# dataset/preprocessed/02a61e1c8dd6698b982d663361c0bb87e0d3522055d16b87bfc4066170e4a8de.json
{
  "hash": "02a61e1c8dd6698b982d663361c0bb87e0d3522055d16b87bfc4066170e4a8de",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  zwj19940929.f3322.org\n  zwj19940929.f3322.org \u2192 104.233.227.23\n\nIP TRAFFIC (1):\n  104.233.227.23:2017 (TCP)\n\nSERVICES OPENED (1):\n  NetRoot\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  services.exe\n    svchost.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (11):\n  - contain obfuscated stackstrings\n  - link function at runtime on Windows\n  - clear Windows event logs\n  - check OS version\n  - start service\n  - create service\n  - persist via Windows service\n  - get memory capacity\n  - parse PE header\n  - link many functions at runtime\n  - reference anti-VM strings\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Process::Create Thread [C0038]\n  - contain loop\n  - link function at runtime on Windows\n  - print debug messages\n  - Discovery::File and Directory Discovery::Log File [E1083.m01]\n  - clear Windows event logs\n  - File System::Writes File [C0052]\n  - File System::Delete File [C0047]\n  - File System::Move File [C0063]\n  - get OS version\n  - Discovery::System Information Discovery [E1082]\n  - get service handle\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - start service\n  - create service\n  - persist via Windows service\n\n--- ARTIFACTS ---\n\nFILES WRITTEN (1):\n  C:\\WINDOWS\\svchost.exe\n\nPERSISTENCE REGISTRY KEYS (6):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\NetRoot\\Group = 5.03\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\NetRoot\\InstallTime = 2019-11-28 06:01\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\NetRoot\\SYSTEM\\CurrentControlSet\\Services\\NetRoot = \ufffd\ufffd\u022b\ufffd\ufffd\ufffd\ufffd\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\NetRoot\\Group = 5.03\u0000\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\NetRoot\\InstallTime = 2019-11-30 14:35\u0000\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\NetRoot\\SYSTEM\\CurrentControlSet\\Services\\NetRoot = E5AE89E585A8E69C8DE58AA100\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\ActiveMovie\\devenum\\Version\n\nMUTEXES CREATED (8):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  eed3bd3a-a1ad-4e99-987b-d7cb3fcfa7f0 - S-1-5-21-1482476501-1645522239-1417001333-500\n  zwj19940929.f3322.org:2017:NetRoot\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc, VirtualProtect\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .rdata       entropy=4.36\n  .data        entropy=6.92\n  hxcysbk      entropy=0.0\n  [overlay]    entropy=0.0  size=14052\n",
  "ground_truth": {
    "hash": "02a61e1c8dd6698b982d663361c0bb87e0d3522055d16b87bfc4066170e4a8de",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8194,
    "malicious_count": 59,
    "popular_threat_category": [
      {
        "value": "trojan",
        "count": 23
      },
      {
        "value": "downloader",
        "count": 9
      }
    ],
    "popular_threat_name": [
      {
        "value": "farfli",
        "count": 8
      },
      {
        "value": "nadrac",
        "count": 8
      },
      {
        "value": "bkdr",
        "count": 2
      }
    ],
    "suggested_label": "trojan.farfli/nadrac",
    "av_family_vote": "win32",
    "has_family_label": true,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1070",
        "signature_description": "clear Windows event logs"
      },
      {
        "id": "T1082",
        "signature_description": "check OS version"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1497",
        "signature_description": "reference anti-VM strings"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1569",
        "signature_description": "create service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "delay execution",
      "contain obfuscated stackstrings",
      "create or open file",
      "reference absolute stream path on Windows",
      "create or open registry key",
      "create thread",
      "contain loop",
      "link function at runtime on Windows",
      "print debug messages",
      "access the Windows event log",
      "clear Windows event logs",
      "write file on Windows",
      "delete file",
      "move file",
      "get OS version",
      "check OS version",
      "get service handle",
      "set registry value",
      "start service",
      "create service",
      "persist via Windows service",
      "copy file",
      "get memory capacity",
      "parse PE header",
      "enumerate PE sections",
      "inspect section memory permissions",
      "link many functions at runtime",
      "rebuild import table",
      "run as service",
      "reference anti-VM strings"
    ],
    "iocs": {
      "ips": [
        "104.233.227.23"
      ],
      "domains": [
        "zwj19940929.f3322.org"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\NetRoot\\Group",
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\NetRoot\\InstallTime",
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\NetRoot\\SYSTEM\\CurrentControlSet\\Services\\NetRoot",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\NetRoot\\Group",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\NetRoot\\InstallTime",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\NetRoot\\SYSTEM\\CurrentControlSet\\Services\\NetRoot"
      ],
      "registry_all_count": 7,
      "files_dropped": [],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "eed3bd3a-a1ad-4e99-987b-d7cb3fcfa7f0 - S-1-5-21-1482476501-1645522239-1417001333-500",
        "zwj19940929.f3322.org:2017:NetRoot"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Jujubox"
  ]
}