# dataset/preprocessed/0439b374b39e66b3657d42585dcb4074e6c483444b029ebf834bc3764f0b0bef.json
{
  "hash": "0439b374b39e66b3657d42585dcb4074e6c483444b029ebf834bc3764f0b0bef",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nIP TRAFFIC (1):\n  <MACHINE_DNS_SERVER>:53 (UDP)\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  0439b374b39e66b3_black cum beast uncut nipples granny  (britney,sonja).mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\0439b374b39e66b3_black cum beast uncut nipples granny  (britney,sonja).mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\0439b374b39e66b3657d42585dcb4074e6c483444b029ebf834bc3764f0b0bef(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\0439b374b39e66b3_black cum beast uncut nipples granny  (britney,sonja).mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african horse xxx [bangbus] stockings .mpg.exe  [sha256: f378a28badf96972...]\n  american animal bukkake [bangbus] young .avi.exe  [sha256: 8affd432e40b43be...]\n  black beast masturbation .rar.exe  [sha256: 16b86ab264efb205...]\n  british action handjob hidden traffic .mpeg.exe  [sha256: 914c939946743dc5...]\n  bukkake lesbian traffic .zip.exe  [sha256: 04f6a97451bbe9ab...]\n  chinese fucking horse [bangbus] shower .mpg.exe  [sha256: 8b1d48c0b9aa04ae...]\n  cum bukkake [bangbus] .mpeg.exe  [sha256: dbe1b2b216d15bef...]\n  cumshot horse sleeping granny .mpg.exe  [sha256: cf874e6e916ae184...]\n  danish blowjob horse several models balls .mpg.exe  [sha256: 81f1497a35e228f0...]\n  danish gay beastiality girls .mpeg.exe  [sha256: 0b581d9434308ef2...]\n  driver.exe  [sha256: 20d437d345496b18...]\n  drwatson32.exe  [sha256: 3842feba6d75264c...]\n  freeftp.exe  [sha256: 6aa6989a28915674...]\n  french xxx gang bang public hole  (Sonja).zip.exe  [sha256: b0576b35d842727b...]\n  game3d.exe  [sha256: fd478373cab1c3b1...]\n  gang bang animal [free] pregnant .mpg.exe  [sha256: ad70663fdeb11077...]\n  gang bang [bangbus] .rar.exe  [sha256: afce2b18c051f390...]\n  hardcore lesbian full movie fishy .mpeg.exe  [sha256: 99b86f3220784065...]\n  horse [bangbus] swallow  (Sonja).mpeg.exe  [sha256: e7c6864b3ca5096d...]\n  indian lesbian gay [bangbus] blondie  (Sylvia).rar.exe  [sha256: 5b293581ddca94c1...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\american animal bukkake [bangbus] young .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\nude [free] glans  (Melissa).rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\norwegian gang bang bukkake girls stockings .rar.exe\n  C:\\Documents and Settings\\All Users\\Templates\\british action handjob hidden traffic .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\nude several models hairy .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\russian cum xxx [free] .rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\malaysia fucking blowjob lesbian  (Ashley,Karin).mpeg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\indian lesbian gay [bangbus] blondie  (Sylvia).rar.exe\n  C:\\Program Files\\Movie Maker\\Shared\\chinese fucking horse [bangbus] shower .mpg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\gang bang animal [free] pregnant .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\bukkake lesbian traffic .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\lingerie beastiality full movie shoes .avi.exe\n  C:\\WINDOWS\\Temp\\cumshot horse sleeping granny .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\hardcore lesbian full movie fishy .mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\spanish animal xxx licking .zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\gang bang [bangbus] .rar.exe\n  C:\\WINDOWS\\ime\\shared\\black beast masturbation .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\horse [bangbus] swallow  (Sonja).mpeg.exe\n  C:\\WINDOWS\\security\\templates\\swedish horse lesbian hot (!) leather .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.880111217498779  size=254718\n",
  "ground_truth": {
    "hash": "0439b374b39e66b3657d42585dcb4074e6c483444b029ebf834bc3764f0b0bef",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [
        "<MACHINE_DNS_SERVER>"
      ],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\danish cum gay uncut mistress .rar.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\trambling girls  (Melissa).mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\indian gang bang horse uncut .avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\bukkake hidden sm .avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\russian fetish beast catfight hole ash  (Samantha).rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\bukkake girls  (Samantha).rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\indian beastiality fucking [milf] hole fishy .mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\black kicking lingerie licking fishy .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "african horse xxx [bangbus] stockings .mpg.exe",
        "american animal bukkake [bangbus] young .avi.exe",
        "black beast masturbation .rar.exe",
        "british action handjob hidden traffic .mpeg.exe",
        "bukkake lesbian traffic .zip.exe",
        "chinese fucking horse [bangbus] shower .mpg.exe",
        "cum bukkake [bangbus] .mpeg.exe",
        "cumshot horse sleeping granny .mpg.exe",
        "danish blowjob horse several models balls .mpg.exe",
        "danish gay beastiality girls .mpeg.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "french xxx gang bang public hole  (Sonja).zip.exe",
        "game3d.exe",
        "gang bang [bangbus] .rar.exe",
        "gang bang animal [free] pregnant .mpg.exe",
        "hardcore lesbian full movie fishy .mpeg.exe",
        "horse [bangbus] swallow  (Sonja).mpeg.exe",
        "indian lesbian gay [bangbus] blondie  (Sylvia).rar.exe",
        "italian beast lingerie full movie femdom .mpeg.exe",
        "lesbian sleeping .zip.exe",
        "lingerie beastiality full movie shoes .avi.exe",
        "malaysia fucking blowjob lesbian  (Ashley,Karin).mpeg.exe",
        "mssrv.exe",
        "norwegian gang bang bukkake girls stockings .rar.exe",
        "notes.txt.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}