# dataset/preprocessed/04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248.json
{
  "hash": "04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (5):\n  ((null)) C:\\04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\04d40de652a0abf1_russian porn horse several models hole castration .rar.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n  \"C:\\Windows\\system32\\rundll32.exe\" C:\\Windows\\system32\\shell32.dll,OpenAs_RunDLL C:\\analyse\\1577503786.5929863_c501ad68-1d9f-426c-86a3-fba27c95ccff\n   \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american cumshot horse [milf] glans circumcision .mpeg.exe  [sha256: 246bdd6a037f7d1a...]\n  beast public pregnant .zip.exe  [sha256: 4d056841c5d9f5b1...]\n  black cumshot lesbian masturbation .rar.exe  [sha256: 1ca13bb20abd5a5a...]\n  brasilian beastiality lingerie big .mpeg.exe  [sha256: 46b5f2684c9b76b3...]\n  brasilian nude blowjob catfight titts sm  (Sylvia).avi.exe  [sha256: f9eabfc74dae865d...]\n  danish beastiality gay catfight feet .rar.exe  [sha256: 2b18a44a32c56b33...]\n  danish porn trambling public  (Karin).mpg.exe  [sha256: 9a44ad1548252cba...]\n  deinstall.exe  [sha256: 30b2e53ac1a15a81...]\n  drwatson32.exe  [sha256: 861a771acbb0c151...]\n  fucking [free] glans .zip.exe  [sha256: 024324193578b575...]\n  gay hidden  (Curtney).zip.exe  [sha256: 5588eaa36ed97c8c...]\n  gay several models penetration .zip.exe  [sha256: a31c1526e50f3dc7...]\n  hardcore catfight .rar.exe  [sha256: be7be1a322943df0...]\n  horse catfight sm .mpeg.exe  [sha256: b8816196c12ced5f...]\n  horse masturbation .mpeg.exe  [sha256: 43cc45ddc7af2e5e...]\n  horse uncut .rar.exe  [sha256: 9d6a45d144d5e58d...]\n  IE9setup.exe  [sha256: 57ab320d15f2eacc...]\n  indian animal beast uncut titts black hairunshaved .avi.exe  [sha256: f491240ed78fece3...]\n  indian fetish lingerie catfight hole 50+ .mpg.exe  [sha256: 12f6b729e18476e1...]\n  indian porn beast several models beautyfull  (Sandy,Karin).rar.exe  [sha256: 1dba0e9a8c9baa38...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\danish beastiality gay catfight feet .rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\italian action horse hot (!) .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\tyrkish animal fucking [bangbus] sm .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\black cumshot lesbian masturbation .rar.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\russian cumshot blowjob several models 50+  (Kathrin,Karin).mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\horse catfight sm .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\indian porn beast several models beautyfull  (Sandy,Karin).rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\indian fetish lingerie catfight hole 50+ .mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\indian animal beast uncut titts black hairunshaved .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\lesbian voyeur hole circumcision .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\lesbian full movie beautyfull .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\horse uncut .rar.exe\n  C:\\WINDOWS\\Temp\\gay hidden  (Curtney).zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\swedish porn bukkake hot (!) titts  (Sandy,Liz).mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\brasilian beastiality lingerie big .mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\tyrkish nude gay masturbation balls  (Anniston,Sylvia).avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\hardcore catfight .rar.exe\n  C:\\WINDOWS\\ime\\shared\\beast public pregnant .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\swedish gang bang blowjob voyeur .avi.exe\n\nPERSISTENCE REGISTRY KEYS (3):\n  HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (4):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n  HKLM/SOFTWARE/Microsoft/Windows/CurrentVersion/Run/mssrv32\n  HKU/S-1-5-21-470376811-3006406624-3672060426-1000/Software/Microsoft/Windows/CurrentVersion/Internet Settings/ZoneMap//U\n  HKU/S-1-5-21-470376811-3006406624-3672060426-1000/Software/Microsoft/Windows/CurrentVersion/Internet Settings/ZoneMap//A\n\nMUTEXES CREATED (14):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  Local\\ZoneAttributeCacheCounterMutex\n  Local\\ZonesCacheCounterMutex\n  Local\\ZonesCounterMutex\n  Local\\ZonesLockedCacheCounterMutex\n\nPROCESSES CREATED (4):\n  C:\\04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248.exe\n  C:\\04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248.exe \n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n  C:\\Windows\\System32\\rundll32.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880280017852783  size=1935305\n",
  "ground_truth": {
    "hash": "04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9437,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 7,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\horse lesbian ejaculation .avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\xxx cum masturbation feet upskirt  (Kathrin).rar.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\african gay lesbian nipples  (Karin).avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\indian gang bang animal full movie gorgeoushorny .mpeg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\lingerie masturbation legs swallow .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\trambling [bangbus] cock stockings  (Kathrin,Sonja).avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\cumshot girls .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\canadian porn hot (!) glans fishy  (Britney,Liz).rar.exe",
        "C:\\Program Files\\dotnet\\shared\\animal lingerie licking granny  (Jenna,Sarah).mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american cumshot horse [milf] glans circumcision .mpeg.exe",
        "beast public pregnant .zip.exe",
        "black cumshot lesbian masturbation .rar.exe",
        "brasilian beastiality lingerie big .mpeg.exe",
        "brasilian nude blowjob catfight titts sm  (Sylvia).avi.exe",
        "c:\\program files\\adobe\\acrobat reader dc\\reader\\idtemplates\\black beast kicking hidden young .mpeg.exe",
        "c:\\program files\\common files\\adobe\\arm\\1.0\\temp\\fucking girls castration  (jade,gina).rar.exe",
        "c:\\program files\\common files\\microsoft shared\\ime14\\shared\\italian nude girls titts bedroom .avi.exe",
        "c:\\program files\\common files\\microsoft shared\\rrloc14\\en-us-shared-cn\\chinese action horse uncut penetration  (anniston).mpg.exe",
        "c:\\program files\\common files\\microsoft shared\\russian bukkake [milf] legs .mpeg.exe",
        "c:\\program files\\dvd maker\\shared\\cumshot girls granny  (samantha,jade).mpg.exe",
        "c:\\windows\\mssrv.exe",
        "danish beastiality gay catfight feet .rar.exe",
        "danish porn trambling public  (Karin).mpg.exe",
        "deinstall.exe",
        "drwatson32.exe",
        "fucking [free] glans .zip.exe",
        "gay hidden  (Curtney).zip.exe",
        "gay several models penetration .zip.exe",
        "hardcore catfight .rar.exe",
        "horse catfight sm .mpeg.exe",
        "horse masturbation .mpeg.exe",
        "horse uncut .rar.exe",
        "indian animal beast uncut titts black hairunshaved .avi.exe",
        "indian fetish lingerie catfight hole 50+ .mpg.exe",
        "indian porn beast several models beautyfull  (Sandy,Karin).rar.exe",
        "italian action horse hot (!) .avi.exe",
        "italian porn sperm lesbian .mpeg.exe",
        "japanese handjob hardcore [bangbus] sm  (Jenna,Liz).mpg.exe",
        "lesbian full movie beautyfull .zip.exe",
        "lesbian masturbation granny .zip.exe",
        "lesbian voyeur hole circumcision .rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian cumshot blowjob several models 50+  (Kathrin,Karin).mpg.exe",
        "screensaver.scr",
        "soccer.exe",
        "start.exe",
        "swedish gang bang blowjob voyeur .avi.exe",
        "swedish porn bukkake hot (!) titts  (Sandy,Liz).mpg.exe",
        "tyrkish animal fucking [bangbus] sm .zip.exe",
        "tyrkish nude gay masturbation balls  (Anniston,Sylvia).avi.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "Local\\ZoneAttributeCacheCounterMutex",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248.exe",
        "C:\\04d40de652a0abf1cc79c56846f155a8acaf98e72076b8e5e4a08439fe27b248.exe ",
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe",
        "C:\\Windows\\System32\\rundll32.exe"
      ]
    }
  },
  "num_sandboxes": 6,
  "sandbox_names": [
    "CAPA",
    "Rising MOVES",
    "Sangfor ZSand",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}