# dataset/preprocessed/066afb06f52abd1ebdb2a12dab3d0783ca7a58c9d3ba9e5a2853719578482351.json
{
  "hash": "066afb06f52abd1ebdb2a12dab3d0783ca7a58c9d3ba9e5a2853719578482351",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  066afb06f52abd1e_black bukkake cumshot catfight hairy  (kathrin).rar.exe\n    C:\\Users\\<USER>\\Downloads\\066afb06f52abd1e_black bukkake cumshot catfight hairy  (kathrin).rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\066afb06f52abd1ebdb2a12dab3d0783ca7a58c9d3ba9e5a2853719578482351(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\066afb06f52abd1e_black bukkake cumshot catfight hairy  (kathrin).rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american porn lingerie girls .avi.exe  [sha256: ed9e053a5c7db5a9...]\n  beast several models glans balls  (Jade).mpeg.exe  [sha256: 7b79810ef5a68bee...]\n  beast voyeur .mpeg.exe  [sha256: 9bc15a82065443bf...]\n  black horse gay public cock shoes .mpeg.exe  [sha256: 81594222e921472a...]\n  blowjob hot (!) hole fishy  (Tatjana).avi.exe  [sha256: c6c39e9d6ba4ebcd...]\n  blowjob [free] .rar.exe  [sha256: 10839b7f4ef0beb3...]\n  brasilian fetish gay catfight .mpeg.exe  [sha256: d7fd062ec567dfa6...]\n  brasilian porn bukkake licking gorgeoushorny .mpeg.exe  [sha256: 958503b1a979375a...]\n  bukkake voyeur cock .mpeg.exe  [sha256: b5699c0dc333fe96...]\n  danish beastiality xxx [milf] hole .rar.exe  [sha256: 3e6409244d7d83ca...]\n  danish horse lingerie lesbian lady .zip.exe  [sha256: db5c5de222085a3c...]\n  defrag64.exe  [sha256: b906e64adebba922...]\n  deinstall.exe  [sha256: a589d1f2c8454841...]\n  driver.exe  [sha256: 8bc858f565ac976e...]\n  fucking sleeping hole girly  (Liz).mpg.exe  [sha256: 4fcc2acc9231a2b8...]\n  hardcore licking girly .rar.exe  [sha256: 4756b0aae0dda7ff...]\n  horse masturbation glans .zip.exe  [sha256: 2acdee78759c9d1f...]\n  horse [free] feet balls .mpeg.exe  [sha256: fa693e14ad5cf0e3...]\n  IE9setup.exe  [sha256: ebc1b794f1af3f85...]\n  indian porn xxx voyeur hole traffic .mpg.exe  [sha256: bd0b02399295869e...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\hardcore licking girly .rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\danish horse lingerie lesbian lady .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\lingerie public .avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\indian porn xxx voyeur hole traffic .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\italian horse fucking catfight redhair .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\japanese action blowjob [milf] 50+ .avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\american porn lingerie girls .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\tyrkish porn blowjob catfight  (Sylvia).rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\italian gang bang lesbian voyeur  (Sylvia).mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\brasilian fetish gay catfight .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\blowjob hot (!) hole fishy  (Tatjana).avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\lingerie several models swallow .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\blowjob [free] .rar.exe\n  C:\\WINDOWS\\Temp\\danish beastiality xxx [milf] hole .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\beast voyeur .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\tyrkish horse sperm licking boots .avi.exe\n  C:\\WINDOWS\\assembly\\temp\\sperm masturbation penetration  (Britney,Liz).avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\horse [free] feet balls .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\russian action fucking sleeping cock ejaculation  (Jade).zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.880297660827637  size=1626569\n",
  "ground_truth": {
    "hash": "066afb06f52abd1ebdb2a12dab3d0783ca7a58c9d3ba9e5a2853719578482351",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9444,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\african fetish licking  (Curtney,Gina).mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\nude masturbation  (Kathrin).zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\hardcore hidden  (Sonja).mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\handjob hidden stockings .avi.exe",
        "C:\\Program Files\\dotnet\\shared\\fetish blowjob full movie .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american porn lingerie girls .avi.exe",
        "beast several models glans balls  (Jade).mpeg.exe",
        "beast voyeur .mpeg.exe",
        "black horse gay public cock shoes .mpeg.exe",
        "blowjob [free] .rar.exe",
        "blowjob hot (!) hole fishy  (Tatjana).avi.exe",
        "brasilian fetish gay catfight .mpeg.exe",
        "brasilian porn bukkake licking gorgeoushorny .mpeg.exe",
        "bukkake voyeur cock .mpeg.exe",
        "danish beastiality xxx [milf] hole .rar.exe",
        "danish horse lingerie lesbian lady .zip.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "fucking sleeping hole girly  (Liz).mpg.exe",
        "hardcore licking girly .rar.exe",
        "horse [free] feet balls .mpeg.exe",
        "horse masturbation glans .zip.exe",
        "indian porn xxx voyeur hole traffic .mpg.exe",
        "italian gang bang lesbian voyeur  (Sylvia).mpeg.exe",
        "italian horse fucking catfight redhair .avi.exe",
        "japanese action blowjob [milf] 50+ .avi.exe",
        "lingerie hidden feet .rar.exe",
        "lingerie public .avi.exe",
        "lingerie several models swallow .rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian action fucking sleeping cock ejaculation  (Jade).zip.exe",
        "screensaver.scr",
        "setup.exe",
        "sperm [bangbus] balls  (Sandy,Janette).avi.exe",
        "sperm masturbation penetration  (Britney,Liz).avi.exe",
        "start.exe",
        "swedish kicking fucking licking hole latex .avi.exe",
        "tyrkish horse sperm licking boots .avi.exe",
        "tyrkish porn blowjob catfight  (Sylvia).rar.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}