# dataset/preprocessed/067b9c39c47ba989f902e68ba483f24a502a2e8a56a0f632dce05a21b8c695d4.json
{
  "hash": "067b9c39c47ba989f902e68ba483f24a502a2e8a56a0f632dce05a21b8c695d4",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  067b9c39c47ba989_brasilian cumshot action big hole fishy .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\067b9c39c47ba989_brasilian cumshot action big hole fishy .mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\067b9c39c47ba989f902e68ba483f24a502a2e8a56a0f632dce05a21b8c695d4(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\067b9c39c47ba989_brasilian cumshot action big hole fishy .mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  beast full movie .rar.exe  [sha256: f69ce278c2035f6e...]\n  beast hidden  (Jade).mpeg.exe  [sha256: d9f03c50ad903d91...]\n  black fetish sperm big redhair  (Gina,Melissa).mpeg.exe  [sha256: 5c5aea77efe2677a...]\n  brasilian gang bang lesbian licking glans hotel  (Sylvia).avi.exe  [sha256: 1049500177652c1c...]\n  bukkake full movie traffic .mpg.exe  [sha256: c9fe6d6c9ac08c2c...]\n  bukkake girls feet bondage  (Sarah).rar.exe  [sha256: 5b636778ade9865f...]\n  bukkake licking .avi.exe  [sha256: 61896ff5bfcaeebc...]\n  danish action gay catfight .mpg.exe  [sha256: ef1824cb3760faf1...]\n  danish nude horse girls hole hotel .zip.exe  [sha256: 65fdb54eea9fb1f7...]\n  defrag64.exe  [sha256: 433ae646c0b629d0...]\n  driver.exe  [sha256: 4ae01ed66bb563ff...]\n  drwatson32.exe  [sha256: 9bd96b34b045b3bf...]\n  freeftp.exe  [sha256: f0993d98bc00dd26...]\n  fucking full movie leather  (Sonja,Liz).avi.exe  [sha256: e6bf7f9578f269d3...]\n  gay full movie 50+ .mpg.exe  [sha256: a6643d924bd78084...]\n  hardcore sleeping  (Janette).avi.exe  [sha256: 1d1ae0c78af1e009...]\n  hardcore [bangbus] titts hairy  (Sarah).zip.exe  [sha256: ec8cc2027985aa79...]\n  hardcore [free] hole granny  (Samantha).mpeg.exe  [sha256: 06600cedd1c96e42...]\n  horse hot (!) shower .mpeg.exe  [sha256: e923143f4021e008...]\n  horse sleeping hole .avi.exe  [sha256: a840c4baa3b1c312...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\tyrkish cum blowjob lesbian hairy .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\italian cum blowjob [milf] feet .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\tyrkish nude xxx big feet stockings .mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\danish nude horse girls hole hotel .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\beast full movie .rar.exe\n  C:\\Documents and Settings\\Default User\\Templates\\xxx catfight lady .rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\brasilian gang bang lesbian licking glans hotel  (Sylvia).avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\hardcore sleeping  (Janette).avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\tyrkish gang bang fucking voyeur upskirt  (Christine,Samantha).rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\gay full movie 50+ .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\russian cum beast [bangbus] .mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\danish action gay catfight .mpg.exe\n  C:\\WINDOWS\\Temp\\hardcore [bangbus] titts hairy  (Sarah).zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\horse hot (!) shower .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\black fetish sperm big redhair  (Gina,Melissa).mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\japanese fetish xxx sleeping feet .mpg.exe\n  C:\\WINDOWS\\assembly\\tmp\\horse sleeping hole .avi.exe\n  C:\\WINDOWS\\ime\\shared\\bukkake girls feet bondage  (Sarah).rar.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\bukkake full movie traffic .mpg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880167484283447  size=514390\n",
  "ground_truth": {
    "hash": "067b9c39c47ba989f902e68ba483f24a502a2e8a56a0f632dce05a21b8c695d4",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\british cumshot [bangbus] legs lady .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\kicking cum sleeping granny  (Sonja,Anniston).mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\british gay licking glans hairy .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\french animal hardcore several models ash circumcision .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\blowjob hidden .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\tyrkish fucking several models .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "beast full movie .rar.exe",
        "beast hidden  (Jade).mpeg.exe",
        "black fetish sperm big redhair  (Gina,Melissa).mpeg.exe",
        "brasilian gang bang lesbian licking glans hotel  (Sylvia).avi.exe",
        "bukkake full movie traffic .mpg.exe",
        "bukkake girls feet bondage  (Sarah).rar.exe",
        "bukkake licking .avi.exe",
        "danish action gay catfight .mpg.exe",
        "danish nude horse girls hole hotel .zip.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking full movie leather  (Sonja,Liz).avi.exe",
        "gay full movie 50+ .mpg.exe",
        "hardcore [bangbus] titts hairy  (Sarah).zip.exe",
        "hardcore [free] hole granny  (Samantha).mpeg.exe",
        "hardcore sleeping  (Janette).avi.exe",
        "horse hot (!) shower .mpeg.exe",
        "horse sleeping hole .avi.exe",
        "italian cum blowjob [milf] feet .zip.exe",
        "italian nude blowjob full movie glans .avi.exe",
        "italian porn sperm lesbian .avi.exe",
        "japanese fetish xxx sleeping feet .mpg.exe",
        "mssrv.exe",
        "rundll64.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}