# dataset/preprocessed/06d331734caace90201c47eb29deb97241b61c6d7bc39aacfa54a3b5f4ca045c.json
{
  "hash": "06d331734caace90201c47eb29deb97241b61c6d7bc39aacfa54a3b5f4ca045c",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nIP TRAFFIC (1):\n  <MACHINE_DNS_SERVER>:53 (UDP)\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\06d331734caace90201c47eb29deb97241b61c6d7bc39aacfa54a3b5f4ca045c(null) [(null)]\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (19):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n  - reference anti-VM strings targeting Xen\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african blowjob fucking masturbation hole sm .zip.exe  [sha256: 93a59bd5a25bcace...]\n  blowjob action full movie  (Tatjana).avi.exe  [sha256: f56b9e9aa0da464a...]\n  brasilian trambling several models .avi.exe  [sha256: 34ca0eace8298d81...]\n  bukkake hot (!)  (Ashley,Jenna).zip.exe  [sha256: e30625bca1b31599...]\n  chinese lesbian sleeping cock sweet  (Sonja).zip.exe  [sha256: 8471138c8e2fdb6e...]\n  danish beast hidden  (Karin).zip.exe  [sha256: 56c69bd2ce74d430...]\n  danish gay hidden latex  (Anniston).mpeg.exe  [sha256: 7ac71618de25b909...]\n  defrag.exe  [sha256: 12bbd83c536d766e...]\n  driver.exe  [sha256: 2ce9bdb5b259355a...]\n  fetish voyeur traffic  (Sandy,Curtney).avi.exe  [sha256: 16000b4a4ea392d7...]\n  french beastiality horse public .mpeg.exe  [sha256: 446472bdf119b0c9...]\n  french fetish cum sleeping young  (Tatjana,Samantha).avi.exe  [sha256: eed17120963c1cf5...]\n  french lingerie lesbian public boobs  (Janette,Samantha).avi.exe  [sha256: 67224c7123bad444...]\n  fucking [free] balls .zip.exe  [sha256: 56604d85525de6d7...]\n  gang bang girls castration  (Kathrin,Sandy).mpeg.exe  [sha256: 5afdbb45aef28f8a...]\n  gay horse voyeur hole traffic .zip.exe  [sha256: cffaa5416969eee7...]\n  hdcleaner.exe  [sha256: bbd8b38a009ef541...]\n  horse licking .mpeg.exe  [sha256: b5e10d4c2fe06541...]\n  horse sleeping 40+ .zip.exe  [sha256: 9a45a4bb8f4c5603...]\n  horse sperm voyeur ash .rar.exe  [sha256: 11dc07042bf3e29a...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\blowjob action full movie  (Tatjana).avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\nude animal [free] .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\danish beast hidden  (Karin).zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\african blowjob fucking masturbation hole sm .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\danish gay hidden latex  (Anniston).mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\french beastiality horse public .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\italian kicking fetish [bangbus] .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\chinese lesbian sleeping cock sweet  (Sonja).zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\japanese gang bang sleeping glans mature  (Janette,Liz).mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\french lingerie lesbian public boobs  (Janette,Samantha).avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\french fetish cum sleeping young  (Tatjana,Samantha).avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\horse sperm voyeur ash .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\gay horse voyeur hole traffic .zip.exe\n  C:\\WINDOWS\\Temp\\kicking girls black hairunshaved  (Karin).rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\spanish cumshot lesbian .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\horse sleeping 40+ .zip.exe\n  C:\\WINDOWS\\assembly\\temp\\indian action action public black hairunshaved .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\fucking [free] balls .zip.exe\n  C:\\WINDOWS\\ime\\shared\\tyrkish handjob hardcore sleeping ash beautyfull  (Samantha,Sarah).mpg.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.8800435066223145  size=1787603\n",
  "ground_truth": {
    "hash": "06d331734caace90201c47eb29deb97241b61c6d7bc39aacfa54a3b5f4ca045c",
    "verdict": "MALICIOUS",
    "detection_rate": 0.942,
    "malicious_count": 65,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1497",
        "signature_description": "reference anti-VM strings targeting Xen"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference anti-VM strings targeting Xen",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [
        "<MACHINE_DNS_SERVER>"
      ],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "IE9setup.exe",
        "african blowjob fucking masturbation hole sm .zip.exe",
        "blowjob action full movie  (Tatjana).avi.exe",
        "brasilian trambling several models .avi.exe",
        "bukkake hot (!)  (Ashley,Jenna).zip.exe",
        "chinese lesbian sleeping cock sweet  (Sonja).zip.exe",
        "danish beast hidden  (Karin).zip.exe",
        "danish gay hidden latex  (Anniston).mpeg.exe",
        "defrag.exe",
        "driver.exe",
        "fetish voyeur traffic  (Sandy,Curtney).avi.exe",
        "french beastiality horse public .mpeg.exe",
        "french fetish cum sleeping young  (Tatjana,Samantha).avi.exe",
        "french lingerie lesbian public boobs  (Janette,Samantha).avi.exe",
        "fucking [free] balls .zip.exe",
        "gang bang girls castration  (Kathrin,Sandy).mpeg.exe",
        "gay horse voyeur hole traffic .zip.exe",
        "hdcleaner.exe",
        "horse licking .mpeg.exe",
        "horse sleeping 40+ .zip.exe",
        "horse sperm voyeur ash .rar.exe",
        "indian action action public black hairunshaved .mpeg.exe",
        "install.exe",
        "italian kicking fetish [bangbus] .avi.exe",
        "japanese action voyeur shower  (Karin).mpeg.exe",
        "japanese cum bukkake several models traffic  (Karin,Samantha).avi.exe",
        "japanese gang bang sleeping glans mature  (Janette,Liz).mpeg.exe",
        "kicking girls black hairunshaved  (Karin).rar.exe",
        "lingerie lingerie licking ash .zip.exe",
        "mssrv.exe",
        "nude animal [free] .mpeg.exe",
        "safetyserver.exe",
        "setup.exe",
        "soccer.exe",
        "spanish cumshot lesbian .mpg.exe",
        "sperm uncut traffic .mpeg.exe",
        "start.exe",
        "summergames.exe",
        "tyrkish blowjob porn hot (!) balls  (Tatjana,Sonja).avi.exe",
        "tyrkish handjob hardcore sleeping ash beautyfull  (Samantha,Sarah).mpg.exe",
        "unformat.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}