# dataset/preprocessed/076fa5172e277d5e218615fc84072b3ea7eb7fe44d48c6787f668a2ffc954dd8.json
{
  "hash": "076fa5172e277d5e218615fc84072b3ea7eb7fe44d48c6787f668a2ffc954dd8",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\076fa5172e277d5e218615fc84072b3ea7eb7fe44d48c6787f668a2ffc954dd8(null) [(null)]\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american beastiality horse hot (!) glans .mpeg.exe  [sha256: a90e711238aacb1d...]\n  american gang bang lingerie lesbian shower .avi.exe  [sha256: cfea5a3a52789dfd...]\n  american kicking hardcore several models cock .mpg.exe  [sha256: 5ed541be46771f4c...]\n  beast catfight glans stockings  (Sarah).zip.exe  [sha256: 33ffc9079d7e1240...]\n  beast uncut feet .mpeg.exe  [sha256: 05707a5a1f18a740...]\n  beast [milf] .avi.exe  [sha256: 793e0191d23d5ee9...]\n  black horse lingerie sleeping feet leather .avi.exe  [sha256: a67fd2f19d908aef...]\n  brasilian animal bukkake several models cock circumcision .mpg.exe  [sha256: f42469026b662094...]\n  brasilian fetish lesbian full movie ejaculation .zip.exe  [sha256: 64c7fdc7a5d3cd6d...]\n  brasilian porn lesbian [free] .zip.exe  [sha256: c00504256468aeff...]\n  brasilian porn sperm voyeur high heels .avi.exe  [sha256: a2abaa2edd60d4b4...]\n  defrag.exe  [sha256: 718acd78dbc1f071...]\n  defrag64.exe  [sha256: 04b584ee1c9bd814...]\n  deinstall.exe  [sha256: 8fd285f43f2fb8b7...]\n  driver.exe  [sha256: d5e004cb98e9042f...]\n  drwatson32.exe  [sha256: a543d70b5135c2e5...]\n  freeftp.exe  [sha256: 80785020cab7fe36...]\n  fucking big hairy .mpeg.exe  [sha256: 11f25b74f88485b7...]\n  gay public cock .avi.exe  [sha256: a430260fb2b5d08b...]\n  hardcore lesbian fishy .avi.exe  [sha256: 1a727866f69b8354...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\beast uncut feet .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\gay public cock .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\tyrkish kicking gay several models mature .rar.exe\n  C:\\Documents and Settings\\All Users\\Templates\\russian beastiality hardcore masturbation titts  (Christine,Samantha).mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\trambling lesbian .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\brasilian fetish lesbian full movie ejaculation .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\beast catfight glans stockings  (Sarah).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\malaysia lingerie uncut hole .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\brasilian porn sperm voyeur high heels .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\swedish beastiality trambling several models glans .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\brasilian porn lesbian [free] .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\beast [milf] .avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\hardcore lesbian fishy .avi.exe\n  C:\\WINDOWS\\Temp\\italian porn hardcore [milf] balls .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\american beastiality horse hot (!) glans .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\trambling catfight titts girly  (Tatjana).mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\fucking big hairy .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\brasilian animal bukkake several models cock circumcision .mpg.exe\n  C:\\WINDOWS\\ime\\shared\\russian kicking hardcore [bangbus] ejaculation .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880361557006836  size=1321856\n",
  "ground_truth": {
    "hash": "076fa5172e277d5e218615fc84072b3ea7eb7fe44d48c6787f668a2ffc954dd8",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "american beastiality horse hot (!) glans .mpeg.exe",
        "american gang bang lingerie lesbian shower .avi.exe",
        "american kicking hardcore several models cock .mpg.exe",
        "beast [milf] .avi.exe",
        "beast catfight glans stockings  (Sarah).zip.exe",
        "beast uncut feet .mpeg.exe",
        "black horse lingerie sleeping feet leather .avi.exe",
        "brasilian animal bukkake several models cock circumcision .mpg.exe",
        "brasilian fetish lesbian full movie ejaculation .zip.exe",
        "brasilian porn lesbian [free] .zip.exe",
        "brasilian porn sperm voyeur high heels .avi.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking big hairy .mpeg.exe",
        "gay public cock .avi.exe",
        "hardcore lesbian fishy .avi.exe",
        "indian beastiality lesbian [milf] ash .avi.exe",
        "indian handjob sperm hot (!) stockings .mpg.exe",
        "italian handjob lesbian big .avi.exe",
        "italian porn hardcore [milf] balls .rar.exe",
        "japanese beastiality bukkake voyeur .rar.exe",
        "japanese beastiality xxx catfight cock traffic .mpg.exe",
        "malaysia lingerie uncut hole .avi.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian beastiality hardcore masturbation titts  (Christine,Samantha).mpeg.exe",
        "russian kicking hardcore [bangbus] ejaculation .zip.exe",
        "setup.exe",
        "start.exe",
        "swedish beastiality trambling several models glans .mpeg.exe",
        "swedish handjob sperm licking high heels .rar.exe",
        "trambling catfight titts girly  (Tatjana).mpeg.exe",
        "trambling lesbian .avi.exe",
        "tyrkish kicking gay several models mature .rar.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}