# dataset/preprocessed/0836ee32d403efcd1a2683d0f101c3f49f8d1d1bd2c57e61841a701e75e8f897.json
{
  "hash": "0836ee32d403efcd1a2683d0f101c3f49f8d1d1bd2c57e61841a701e75e8f897",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  0836ee32d403efcd_cum girls feet upskirt  (sylvia).mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\0836ee32d403efcd_cum girls feet upskirt  (sylvia).mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\0836ee32d403efcd1a2683d0f101c3f49f8d1d1bd2c57e61841a701e75e8f897(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\0836ee32d403efcd_cum girls feet upskirt  (sylvia).mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american fetish gay hot (!) glans .mpg.exe  [sha256: 9583509a01b9ab53...]\n  american fetish xxx catfight hotel .mpg.exe  [sha256: a6855c1b63486bc4...]\n  american nude trambling voyeur .mpg.exe  [sha256: 2bd5941082e98b5e...]\n  american porn lingerie girls sweet .rar.exe  [sha256: 371664213b54c482...]\n  black action hardcore hot (!) glans .zip.exe  [sha256: a8799c175b708f38...]\n  black beastiality bukkake licking .mpeg.exe  [sha256: 3747d709561f6217...]\n  brasilian handjob fucking full movie feet  (Anniston,Jade).mpeg.exe  [sha256: 5c12269714462a5c...]\n  danish beastiality lesbian girls .mpeg.exe  [sha256: 6a5b5d76a25c290f...]\n  danish porn lingerie catfight hole fishy  (Samantha).mpg.exe  [sha256: 19fee8377c0095e4...]\n  defrag.exe  [sha256: 4537f0db93838dcc...]\n  deinstall.exe  [sha256: 5bd704edc0dcaae4...]\n  driver.exe  [sha256: b50ff75f883ec055...]\n  freeftp.exe  [sha256: d918ad181b2c8bf1...]\n  fucking sleeping .zip.exe  [sha256: 01d69d95d34f9989...]\n  hardcore girls 50+ .mpg.exe  [sha256: fc5c3791faac5b30...]\n  hardcore lesbian  (Sylvia).avi.exe  [sha256: 3adfe4dada3afa22...]\n  hardcore [bangbus] feet pregnant .rar.exe  [sha256: 150b9c976a5c64a9...]\n  IE9setup.exe  [sha256: 6211c2e3bf0592a1...]\n  indian action horse [bangbus] feet .mpg.exe  [sha256: c164cc100b71835f...]\n  indian cum horse [milf] mistress .mpg.exe  [sha256: 363467f13a0e3808...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\sperm girls glans beautyfull .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\indian cum horse [milf] mistress .mpg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\swedish gang bang trambling uncut sweet .mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\hardcore girls 50+ .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\indian action horse [bangbus] feet .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\sperm [milf] mature .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\american porn lingerie girls sweet .rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\xxx sleeping  (Tatjana).mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\american nude trambling voyeur .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\indian nude gay lesbian hairy .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\swedish cumshot hardcore hidden feet .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\danish porn lingerie catfight hole fishy  (Samantha).mpg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\fucking sleeping .zip.exe\n  C:\\WINDOWS\\Temp\\black beastiality bukkake licking .mpeg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\lingerie hidden glans bondage .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\russian cumshot lesbian licking bondage  (Kathrin,Samantha).mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\hardcore lesbian  (Sylvia).avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\sperm girls sm .avi.exe\n  C:\\WINDOWS\\ime\\shared\\hardcore [bangbus] feet pregnant .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.64\n  [overlay]    entropy=7.880162239074707  size=481354\n",
  "ground_truth": {
    "hash": "0836ee32d403efcd1a2683d0f101c3f49f8d1d1bd2c57e61841a701e75e8f897",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\lesbian hot (!)  (Sylvia).mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\fucking public 40+ .rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\american gang bang hardcore girls .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\tyrkish nude lingerie [bangbus] hotel .avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\japanese cumshot lesbian [milf] upskirt .mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\danish kicking bukkake [bangbus] .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american fetish gay hot (!) glans .mpg.exe",
        "american fetish xxx catfight hotel .mpg.exe",
        "american nude trambling voyeur .mpg.exe",
        "american porn lingerie girls sweet .rar.exe",
        "black action hardcore hot (!) glans .zip.exe",
        "black beastiality bukkake licking .mpeg.exe",
        "brasilian handjob fucking full movie feet  (Anniston,Jade).mpeg.exe",
        "danish beastiality lesbian girls .mpeg.exe",
        "danish porn lingerie catfight hole fishy  (Samantha).mpg.exe",
        "defrag.exe",
        "deinstall.exe",
        "driver.exe",
        "freeftp.exe",
        "fucking sleeping .zip.exe",
        "hardcore [bangbus] feet pregnant .rar.exe",
        "hardcore girls 50+ .mpg.exe",
        "hardcore lesbian  (Sylvia).avi.exe",
        "indian action horse [bangbus] feet .mpg.exe",
        "indian cum horse [milf] mistress .mpg.exe",
        "indian nude gay lesbian hairy .mpeg.exe",
        "italian beastiality xxx licking glans .mpeg.exe",
        "italian porn lingerie big stockings  (Kathrin,Sarah).mpg.exe",
        "lesbian big fishy .mpg.exe",
        "lesbian hidden .avi.exe",
        "lingerie hidden glans bondage .mpeg.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian cumshot lesbian licking bondage  (Kathrin,Samantha).mpg.exe",
        "screensaver.scr",
        "soccer.exe",
        "sperm [milf] mature .zip.exe",
        "sperm girls glans beautyfull .mpeg.exe",
        "sperm girls sm .avi.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}