# dataset/preprocessed/083924bd9101c0dcb0665a44c5db34a8b25392c6bee083f8f5f65456ac5d67cf.json
{
  "hash": "083924bd9101c0dcb0665a44c5db34a8b25392c6bee083f8f5f65456ac5d67cf",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  083924bd9101c0dc_italian fetish sperm girls hole young .mpg.exe\n    C:\\Users\\<USER>\\Downloads\\083924bd9101c0dc_italian fetish sperm girls hole young .mpg.exe\n    C:\\Users\\<USER>\\Downloads\\083924bd9101c0dc_italian fetish sperm girls hole young .mpg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\083924bd9101c0dcb0665a44c5db34a8b25392c6bee083f8f5f65456ac5d67cf(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\083924bd9101c0dc_italian fetish sperm girls hole young .mpg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american beastiality sperm licking  (Sylvia).avi.exe  [sha256: fd4b331a38589355...]\n  american nude fucking several models glans sweet .mpeg.exe  [sha256: b0f295b7d2939709...]\n  black cum xxx public .zip.exe  [sha256: 9d3b3b29439b7511...]\n  blowjob hot (!)  (Sarah).avi.exe  [sha256: 29872ca9ded83dc5...]\n  bukkake full movie upskirt  (Britney,Sarah).avi.exe  [sha256: 4a1b8cf788e3e0f4...]\n  bukkake hot (!) cock penetration .avi.exe  [sha256: d38b757e16694617...]\n  defrag.exe  [sha256: a18d442990276b1a...]\n  defrag64.exe  [sha256: 5f7971f68bbdc581...]\n  driver.exe  [sha256: 5c96e918eac52e05...]\n  drwatson32.exe  [sha256: fc631980f3d6bdb4...]\n  fucking catfight bedroom .avi.exe  [sha256: 60e217ad675dc369...]\n  fucking [free] titts upskirt .avi.exe  [sha256: 7d9a4dc62765de29...]\n  gay girls cock high heels .avi.exe  [sha256: c7ef6b417cd8681d...]\n  gay [milf] .mpg.exe  [sha256: 093a4b6dd86000a8...]\n  horse catfight blondie .mpeg.exe  [sha256: 17e0b1b56ac2d79b...]\n  IE9setup.exe  [sha256: 858afbc51daaf4ea...]\n  indian gang bang sperm lesbian hole traffic .avi.exe  [sha256: bbdf7cb624c2274f...]\n  indian gang bang xxx [milf] shower  (Anniston,Sarah).zip.exe  [sha256: 72744d8a08b983b9...]\n  italian action hardcore [bangbus] granny .mpeg.exe  [sha256: 9b47cba5f00bf80f...]\n  italian animal lingerie [bangbus] cock fishy .zip.exe  [sha256: 8920adeaf0751f58...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\swedish horse lingerie [free] sweet .zip.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\tyrkish beastiality beast full movie fishy .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\italian kicking sperm hot (!) .mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\italian action hardcore [bangbus] granny .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\japanese horse xxx public .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\blowjob hot (!)  (Sarah).avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\indian gang bang xxx [milf] shower  (Anniston,Sarah).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\italian animal lingerie [bangbus] cock fishy .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\gay girls cock high heels .avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\black cum xxx public .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\sperm [milf] glans traffic .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\tyrkish handjob hardcore lesbian .mpeg.exe\n  C:\\WINDOWS\\Temp\\fucking [free] titts upskirt .avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\swedish animal trambling hot (!) cock .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\spanish trambling [free] titts .avi.exe\n  C:\\WINDOWS\\assembly\\temp\\fucking catfight bedroom .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\trambling sleeping titts ash  (Liz).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\indian gang bang sperm lesbian hole traffic .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\swedish cum lesbian big titts .mpeg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880451202392578  size=1989187\n",
  "ground_truth": {
    "hash": "083924bd9101c0dcb0665a44c5db34a8b25392c6bee083f8f5f65456ac5d67cf",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9429,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\spanish gay trambling lesbian young .mpeg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\spanish beast sperm licking leather  (Kathrin).mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\chinese gang bang beastiality licking fishy  (Sonja).avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\kicking [bangbus] cock .avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\american action cum [free] 50+ .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\brasilian kicking fetish masturbation ash  (Gina).mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\spanish gang bang catfight nipples ejaculation .mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\british bukkake xxx voyeur legs .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\american porn girls .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american beastiality sperm licking  (Sylvia).avi.exe",
        "american nude fucking several models glans sweet .mpeg.exe",
        "black cum xxx public .zip.exe",
        "blowjob hot (!)  (Sarah).avi.exe",
        "bukkake full movie upskirt  (Britney,Sarah).avi.exe",
        "bukkake hot (!) cock penetration .avi.exe",
        "defrag.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking [free] titts upskirt .avi.exe",
        "fucking catfight bedroom .avi.exe",
        "gay [milf] .mpg.exe",
        "gay girls cock high heels .avi.exe",
        "horse catfight blondie .mpeg.exe",
        "indian gang bang sperm lesbian hole traffic .avi.exe",
        "indian gang bang xxx [milf] shower  (Anniston,Sarah).zip.exe",
        "italian action hardcore [bangbus] granny .mpeg.exe",
        "italian animal lingerie [bangbus] cock fishy .zip.exe",
        "italian kicking sperm hot (!) .mpg.exe",
        "japanese horse xxx public .zip.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "setup.exe",
        "soccer.exe",
        "spanish trambling [free] titts .avi.exe",
        "sperm [milf] glans traffic .rar.exe",
        "start.exe",
        "swedish animal trambling hot (!) cock .mpg.exe",
        "swedish cum lesbian big titts .mpeg.exe",
        "swedish horse lingerie [free] sweet .zip.exe",
        "trambling full movie upskirt .avi.exe",
        "trambling sleeping titts ash  (Liz).mpg.exe",
        "tyrkish beastiality beast full movie fishy .mpeg.exe",
        "tyrkish beastiality gay voyeur sweet  (Jenna,Sylvia).mpg.exe",
        "tyrkish handjob hardcore lesbian .mpeg.exe",
        "tyrkish nude blowjob [free] hole hairy .mpg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}