# dataset/preprocessed/08c63271d1359aca7bdb5002d61ca73c5f3110462bfa577c458e16f97c45de96.json
{
  "hash": "08c63271d1359aca7bdb5002d61ca73c5f3110462bfa577c458e16f97c45de96",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    whkjqr.exe\n\nCOMMAND EXECUTIONS (1):\n  (open) C:\\08c63271d1359aca7bdb5002d61ca73c5f3110462bfa577c458e16f97c4(null) [(null)]\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - contain obfuscated stackstrings\n  - get startup folder\n  - set file attributes\n  - get common file path\n  - accept command line arguments\n\nCAPA RULE MATCHES (MBC):\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - get startup folder\n  - Process::Create Process [C0017]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - contain loop\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Delete File [C0047]\n  - File System::Create Directory [C0046]\n  - Process::Allocate Thread Local Storage [C0040]\n  - get thread local storage value\n  - Process::Set Thread Local Storage Value [C0041]\n  - File System::Read File [C0051]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - resolve function by parsing PE exports\n  - contain a resource (.rsrc) section\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  AUTOEXEC.BAT .exe  [sha256: 90e8e7c9abcaa434...]\n  boot.ini .exe  [sha256: c13af277432c007f...]\n  CONFIG.SYS .exe  [sha256: 0cceffc2eaadcaae...]\n  cuckoo .exe  [sha256: ee38db617b065a8e...]\n  Documents and Settings .exe  [sha256: 0b26a0de0c0cacbe...]\n  eula.1028.txt .exe  [sha256: 8b2f3e63711a14fe...]\n  eula.1031.txt .exe  [sha256: d3930e210b737e4d...]\n  eula.1033.txt .exe  [sha256: 0e609d4449321c6d...]\n  eula.1036.txt .exe  [sha256: 99f98ec7de8cf16e...]\n  eula.1040.txt .exe  [sha256: 9538fae3868c1227...]\n  eula.1041.txt .exe  [sha256: ab67249c6e451ef5...]\n  eula.1042.txt .exe  [sha256: 3d0ef7673250605e...]\n  eula.2052.txt .exe  [sha256: 443b138860716fa0...]\n  eula.3082.txt .exe  [sha256: c5688d843ecdb165...]\n  globdata.ini .exe  [sha256: e3d83ea7c7c3150d...]\n  install.exe .exe  [sha256: 023c8cc50c0c5bd8...]\n  install.ini .exe  [sha256: 31770b8a984d9039...]\n  install.res.1028.dll .exe  [sha256: 41481b80a5323725...]\n  install.res.1031.dll .exe  [sha256: 858014faa9731ac2...]\n  install.res.1033.dll .exe  [sha256: 1f160e88deee5505...]\n\nFILES WRITTEN (30):\n  C:\\0 .exe\n  C:\\AUTOEXEC.BAT .exe\n  C:\\CONFIG.SYS .exe\n  C:\\DiskD .exe\n  C:\\DiskX .exe\n  C:\\DllLoader.exe .exe\n  C:\\Documents and Settings .exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Saaaalamm\\Mira.h\n  C:\\Documents and Settings\\All Users\\Application Data\\strme.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\whkjqr.exe\n  C:\\IO.SYS .exe\n  C:\\MSDOS.SYS .exe\n  C:\\Miraa-newfile .exe\n  C:\\Miraa-samplefile .exe\n  C:\\NTDETECT.COM .exe\n  C:\\Program Files .exe\n  C:\\Python27 .exe\n  C:\\boot.ini .exe\n  C:\\bootfont.bin .exe\n  C:\\eula.2052.txt .exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System = C:\\Documents and Settings\\All Users\\Application Data\\whkjqr.exe\u0000\n\nMUTEXES CREATED (7):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\All Users\\Application Data\\strme.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA\n\nCOMPILER AND PACKER:\n  [Compiler] MinGW\n  [Linker] GNU linker ld (GNU Binutils)\n  Packers: {'PEiD': 'MingWin32 v?.? (h)'}\n\nPE SECTIONS (entropy):\n  .text        entropy=6.08\n  .data        entropy=1.22\n  .rdata       entropy=5.01\n  .bss         entropy=0.0\n  .idata       entropy=4.29\n  .rsrc        entropy=7.01\n  .NewSec      entropy=6.99\n  [overlay]    entropy=6.783712387084961  size=310694\n",
  "ground_truth": {
    "hash": "08c63271d1359aca7bdb5002d61ca73c5f3110462bfa577c458e16f97c45de96",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8732,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "get startup folder"
      }
    ],
    "capa_sig_names": [
      "contain obfuscated stackstrings",
      "get startup folder",
      "create process on Windows",
      "set file attributes",
      "delay execution",
      "create or open file",
      "contain loop",
      "get common file path",
      "delete file",
      "create directory",
      "allocate thread local storage",
      "get thread local storage value",
      "set thread local storage value",
      "read file on Windows",
      "accept command line arguments",
      "resolve function by parsing PE exports",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "AUTOEXEC.BAT .exe",
        "CONFIG.SYS .exe",
        "Documents and Settings .exe",
        "IO.SYS .exe",
        "MSDOS.SYS .exe",
        "Mira.h",
        "NTDETECT.COM .exe",
        "Program Files .exe",
        "Python27 .exe",
        "RECYCLER .exe",
        "System Volume Information .exe",
        "VC_RED.MSI .exe",
        "VC_RED.cab .exe",
        "WINDOWS .exe",
        "boot.ini .exe",
        "cuckoo .exe",
        "eula.1028.txt .exe",
        "eula.1031.txt .exe",
        "eula.1033.txt .exe",
        "eula.1036.txt .exe",
        "eula.1040.txt .exe",
        "eula.1041.txt .exe",
        "eula.1042.txt .exe",
        "eula.2052.txt .exe",
        "eula.3082.txt .exe",
        "globdata.ini .exe",
        "install.exe .exe",
        "install.ini .exe",
        "install.res.1028.dll .exe",
        "install.res.1031.dll .exe",
        "install.res.1033.dll .exe",
        "install.res.1036.dll .exe",
        "install.res.1040.dll .exe",
        "install.res.1041.dll .exe",
        "install.res.1042.dll .exe",
        "install.res.2052.dll .exe",
        "install.res.3082.dll .exe",
        "ntldr .exe",
        "pagefile.sys .exe",
        "strme.exe",
        "vcredist.bmp .exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\Documents and Settings\\All Users\\Application Data\\strme.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}