# dataset/preprocessed/08e38e20403384161df6911a8cb8f7a3daf25a7bb836193ae6254270bb29f7f6.json
{
  "hash": "08e38e20403384161df6911a8cb8f7a3daf25a7bb836193ae6254270bb29f7f6",
  "report_text": "--- NETWORK ACTIVITY ---\n\nSERVICES OPENED (1):\n  W32TIME\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  c:\\l.bat .exe\n    C:\\l.bat\n      w32tm /resync /rediscover\n        w32tm /resync /rediscover\n      net start w32time\n        C:\\Windows\\system32\\net1 start w32time\n      w32tm /config /reliable:yes\n        w32tm /config /reliable:yes\n      w32tm /config /syncfromflags:manual /manualpeerlist:10.0.2.2\n        w32tm /config /syncfromflags:manual /manualpeerlist:10.0.2.2\n      net stop w32time\n        C:\\Windows\\system32\\net1 stop w32time\n      cscript /B C:\\Windows\\\\system32\\\\slmgr.vbs /rearm\n      cscript /B C:\\Windows\\\\system32\\\\slmgr.vbs /cpky\n      cscript /B C:\\Windows\\\\system32\\\\slmgr.vbs /upk\n\nCOMMAND EXECUTIONS (2):\n  (open) C:\\08e38e20403384161df6911a8cb8f7a3daf25a7bb836193ae6254270bb2(null) [(null)]\n  c:\\l.bat\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - contain obfuscated stackstrings\n  - get startup folder\n  - set file attributes\n  - get common file path\n  - accept command line arguments\n\nCAPA RULE MATCHES (MBC):\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - get startup folder\n  - Process::Create Process [C0017]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - contain loop\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Delete File [C0047]\n  - File System::Create Directory [C0046]\n  - Process::Allocate Thread Local Storage [C0040]\n  - get thread local storage value\n  - Process::Set Thread Local Storage Value [C0041]\n  - File System::Read File [C0051]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - resolve function by parsing PE exports\n  - contain a resource (.rsrc) section\n\n--- ARTIFACTS ---\n\nFILES DROPPED (4):\n  gmlrqw.exe  [sha256: 23ca6d22408f89a6...]\n  Mira.h  [sha256: 93ba5be8e47ad44f...]\n  C:\\ProgramData\\iacih.exe  [sha256: 69170dbbfcd00db0...]\n  C:\\ProgramData\\Saaaalamm\\Mira.h  [sha256: 93ba5be8e47ad44f...]\n\nFILES WRITTEN (2):\n  C:\\Documents and Settings\\All Users\\Application Data\\Saaaalamm\\Mira.h\n  C:\\Documents and Settings\\All Users\\Application Data\\gmlrqw.exe\n\nPERSISTENCE REGISTRY KEYS (3):\n  Control Panel\\Desktop\\LanguageConfiguration\\System\\CurrentControlSet\\Services\\W32Time\\Config\\AnnounceFlags = 5\n  Control Panel\\Desktop\\LanguageConfiguration\\System\\CurrentControlSet\\Services\\W32Time\\Parameters\\NtpServer = 1\n  Control Panel\\Desktop\\LanguageConfiguration\\System\\CurrentControlSet\\Services\\W32Time\\Parameters\\Type = N\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\All Users\\Application Data\\gmlrqw.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA\n\nCOMPILER AND PACKER:\n  [Compiler] MinGW\n  [Linker] GNU linker ld (GNU Binutils)\n  Packers: {'PEiD': 'MingWin32 v?.? (h)'}\n\nPE SECTIONS (entropy):\n  .text        entropy=6.08\n  .data        entropy=0.71\n  .rdata       entropy=5.01\n  .bss         entropy=0.0\n  .idata       entropy=4.29\n  .rsrc        entropy=6.05\n  [overlay]    entropy=7.12695837020874  size=102662\n",
  "ground_truth": {
    "hash": "08e38e20403384161df6911a8cb8f7a3daf25a7bb836193ae6254270bb29f7f6",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9028,
    "malicious_count": 65,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "get startup folder"
      }
    ],
    "capa_sig_names": [
      "contain obfuscated stackstrings",
      "get startup folder",
      "create process on Windows",
      "set file attributes",
      "delay execution",
      "create or open file",
      "contain loop",
      "get common file path",
      "delete file",
      "create directory",
      "allocate thread local storage",
      "get thread local storage value",
      "set thread local storage value",
      "read file on Windows",
      "accept command line arguments",
      "resolve function by parsing PE exports",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "Control Panel\\Desktop\\LanguageConfiguration\\System\\CurrentControlSet\\Services\\W32Time\\Config\\AnnounceFlags",
        "Control Panel\\Desktop\\LanguageConfiguration\\System\\CurrentControlSet\\Services\\W32Time\\Parameters\\NtpServer",
        "Control Panel\\Desktop\\LanguageConfiguration\\System\\CurrentControlSet\\Services\\W32Time\\Parameters\\Type"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\ProgramData\\Saaaalamm\\Mira.h",
        "C:\\ProgramData\\iacih.exe",
        "Mira.h",
        "gmlrqw.exe"
      ],
      "mutexes": [],
      "processes": [
        "C:\\Documents and Settings\\All Users\\Application Data\\gmlrqw.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}