# dataset/preprocessed/094f0b876c5e87b16b540473dc4ed4a2f21725ef57e9bee9247b336789c9d166.json
{
  "hash": "094f0b876c5e87b16b540473dc4ed4a2f21725ef57e9bee9247b336789c9d166",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  094f0b876c5e87b1_xxx several models feet .avi.exe\n    C:\\Users\\<USER>\\Downloads\\094f0b876c5e87b1_xxx several models feet .avi.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\094f0b876c5e87b16b540473dc4ed4a2f21725ef57e9bee9247b336789c9d166(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\094f0b876c5e87b1_xxx several models feet .avi.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african bukkake big lady  (Gina,Tatjana).avi.exe  [sha256: e2ea5b8073037bf1...]\n  american cum blowjob hot (!) sweet  (Ashley,Liz).mpeg.exe  [sha256: 6143265f184b7d94...]\n  american cumshot lingerie [free] cock .rar.exe  [sha256: 35e4fdcb064dcbcd...]\n  american fetish beast [milf] hole hairy .avi.exe  [sha256: b52532bf3a5e1371...]\n  blowjob several models hole .zip.exe  [sha256: 1558e8898956f562...]\n  brasilian animal xxx several models cock .zip.exe  [sha256: 61b196e36cc1b2f3...]\n  brasilian cum trambling girls titts shoes .rar.exe  [sha256: 727f8ae2ba17f317...]\n  brasilian porn horse [bangbus]  (Jade).rar.exe  [sha256: 1c4b9b0ca86220cc...]\n  danish nude xxx licking high heels .zip.exe  [sha256: 9ab9a9e376e0da19...]\n  defrag64.exe  [sha256: 3089776225b477df...]\n  deinstall.exe  [sha256: d5ac242a4cadb4a4...]\n  drwatson32.exe  [sha256: 5d95090b9a708369...]\n  freeftp.exe  [sha256: 98529fc9e99b667b...]\n  fucking licking cock stockings .rar.exe  [sha256: 5220ba2c774c1eb4...]\n  gay public titts .mpg.exe  [sha256: 26337b46a8c0b223...]\n  hardcore sleeping latex .mpeg.exe  [sha256: 26d47032916cc696...]\n  horse hot (!) feet  (Kathrin,Liz).avi.exe  [sha256: 834f3da1db9dbc6b...]\n  indian cum beast sleeping titts  (Ashley,Karin).rar.exe  [sha256: 6bdf38e725a1e7ab...]\n  italian kicking xxx [milf] cock mature .rar.exe  [sha256: cf939a76a6685c8d...]\n  japanese cumshot lingerie several models .mpg.exe  [sha256: 7ab88a2886b2d6c1...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\brasilian cum trambling girls titts shoes .rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\brasilian animal xxx several models cock .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\american fetish beast [milf] hole hairy .avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\horse hot (!) feet  (Kathrin,Liz).avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\tyrkish animal bukkake big 40+ .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\danish nude xxx licking high heels .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\lingerie [free] glans balls .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\swedish horse lesbian hidden titts hairy .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\trambling hidden feet  (Sonja,Samantha).mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\japanese cumshot lingerie several models .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\japanese gang bang hardcore hidden high heels .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\brasilian porn horse [bangbus]  (Jade).rar.exe\n  C:\\WINDOWS\\Temp\\sperm girls hole  (Christine,Tatjana).rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\american cum blowjob hot (!) sweet  (Ashley,Liz).mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\lingerie girls blondie .avi.exe\n  C:\\WINDOWS\\assembly\\temp\\italian kicking xxx [milf] cock mature .rar.exe\n  C:\\WINDOWS\\assembly\\tmp\\japanese gang bang lingerie lesbian .mpg.exe\n  C:\\WINDOWS\\ime\\shared\\hardcore sleeping latex .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\swedish fetish lesbian masturbation .mpg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.64\n  [overlay]    entropy=7.880500793457031  size=555330\n",
  "ground_truth": {
    "hash": "094f0b876c5e87b16b540473dc4ed4a2f21725ef57e9bee9247b336789c9d166",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\american sperm bukkake licking penetration .rar.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\canadian lingerie lesbian  (Melissa).rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\british beast hidden pregnant .zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\handjob cumshot hot (!) pregnant .avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\russian trambling hidden upskirt  (Kathrin,Melissa).zip.exe",
        "C:\\Program Files\\dotnet\\shared\\danish lingerie fetish sleeping ash upskirt .mpeg.exe",
        "C:\\Windows\\mssrv.exe",
        "african bukkake big lady  (Gina,Tatjana).avi.exe",
        "american cum blowjob hot (!) sweet  (Ashley,Liz).mpeg.exe",
        "american cumshot lingerie [free] cock .rar.exe",
        "american fetish beast [milf] hole hairy .avi.exe",
        "blowjob several models hole .zip.exe",
        "brasilian animal xxx several models cock .zip.exe",
        "brasilian cum trambling girls titts shoes .rar.exe",
        "brasilian porn horse [bangbus]  (Jade).rar.exe",
        "danish nude xxx licking high heels .zip.exe",
        "defrag64.exe",
        "deinstall.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking licking cock stockings .rar.exe",
        "gay public titts .mpg.exe",
        "hardcore sleeping latex .mpeg.exe",
        "horse hot (!) feet  (Kathrin,Liz).avi.exe",
        "indian cum beast sleeping titts  (Ashley,Karin).rar.exe",
        "italian kicking xxx [milf] cock mature .rar.exe",
        "japanese cumshot lingerie several models .mpg.exe",
        "japanese gang bang hardcore hidden high heels .rar.exe",
        "japanese gang bang lingerie lesbian .mpg.exe",
        "lingerie [free] glans balls .avi.exe",
        "lingerie girls blondie .avi.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "screensaver.scr",
        "setup.exe",
        "sperm big feet .rar.exe",
        "sperm girls hole  (Christine,Tatjana).rar.exe",
        "start.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}