# dataset/preprocessed/095b7582f0f1b0886ea6b35718d4eb93476b4da3d471642ce8362d18fa79e8e9.json
{
  "hash": "095b7582f0f1b0886ea6b35718d4eb93476b4da3d471642ce8362d18fa79e8e9",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  tyrkish nude blowjob [free] hole hairy .mpg.exe\n    C:\\Users\\<USER>\\Downloads\\tyrkish nude blowjob [free] hole hairy .mpg.exe\n    C:\\Users\\<USER>\\Downloads\\tyrkish nude blowjob [free] hole hairy .mpg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\095b7582f0f1b0886ea6b35718d4eb93476b4da3d471642ce8362d18fa79e8e9(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\tyrkish nude blowjob [free] hole hairy .mpg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american action bukkake full movie castration  (Gina,Jade).mpeg.exe  [sha256: cba0b2acf90d6b9a...]\n  american beastiality horse masturbation titts femdom  (Tatjana).mpg.exe  [sha256: 44c559e8c01a94f0...]\n  american nude hardcore girls glans penetration  (Janette).rar.exe  [sha256: 608377ec5ed30b0e...]\n  beast uncut leather .mpg.exe  [sha256: d0a30e04800291c2...]\n  black handjob beast big wifey .avi.exe  [sha256: c27eb8d6bce193b7...]\n  black handjob hardcore sleeping feet femdom  (Karin).mpeg.exe  [sha256: d37826bd19c462d7...]\n  brasilian action beast uncut feet  (Sonja,Tatjana).rar.exe  [sha256: 36bac27ce3adf4a4...]\n  brasilian horse gay uncut  (Jade).rar.exe  [sha256: 177959f8f6b2c5e8...]\n  brasilian porn xxx [bangbus] bondage .mpg.exe  [sha256: 53230c52161e7595...]\n  bukkake lesbian  (Curtney).mpeg.exe  [sha256: 9461728cb1dde55b...]\n  danish cum sperm licking .mpeg.exe  [sha256: c259fd6baf562228...]\n  danish nude trambling voyeur .avi.exe  [sha256: 5b0678efe3a348ec...]\n  defrag64.exe  [sha256: c03414a981eabde4...]\n  deinstall.exe  [sha256: 29125f59881cdfdb...]\n  driver.exe  [sha256: c0fbb8ecbffd4232...]\n  drwatson32.exe  [sha256: 67361fbbbe793966...]\n  freeftp.exe  [sha256: 34623fdf3356588c...]\n  fucking several models cock .rar.exe  [sha256: 3548fcf08b11b99f...]\n  gay [milf] glans .mpg.exe  [sha256: 9060b5011970aecf...]\n  horse hidden feet bondage .rar.exe  [sha256: ec674dc342d56356...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\trambling lesbian titts wifey .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\bukkake lesbian  (Curtney).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\russian kicking xxx full movie beautyfull .mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\american nude hardcore girls glans penetration  (Janette).rar.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\danish nude trambling voyeur .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\black handjob hardcore sleeping feet femdom  (Karin).mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\tyrkish porn horse uncut cock .zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\american beastiality horse masturbation titts femdom  (Tatjana).mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\brasilian horse gay uncut  (Jade).rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\beast uncut leather .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\tyrkish action xxx licking redhair .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\lingerie licking pregnant .avi.exe\n  C:\\WINDOWS\\Temp\\russian cum beast lesbian pregnant .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\gay [milf] glans .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\american action bukkake full movie castration  (Gina,Jade).mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\brasilian action beast uncut feet  (Sonja,Tatjana).rar.exe\n  C:\\WINDOWS\\assembly\\tmp\\horse hidden feet bondage .rar.exe\n  C:\\WINDOWS\\ime\\shared\\tyrkish cum horse girls feet stockings .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\fucking several models cock .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880133152008057  size=425809\n",
  "ground_truth": {
    "hash": "095b7582f0f1b0886ea6b35718d4eb93476b4da3d471642ce8362d18fa79e8e9",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\fucking full movie  (Karin).rar.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\swedish cum xxx big  (Jade).mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\russian beastiality xxx several models bedroom .mpeg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\italian animal bukkake big  (Samantha).rar.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\swedish beastiality bukkake masturbation feet  (Ashley,Jade).avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\fucking licking high heels  (Sandy,Tatjana).mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\blowjob sleeping  (Curtney).mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\bukkake full movie glans black hairunshaved .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\american beastiality beast several models glans .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "american action bukkake full movie castration  (Gina,Jade).mpeg.exe",
        "american beastiality horse masturbation titts femdom  (Tatjana).mpg.exe",
        "american nude hardcore girls glans penetration  (Janette).rar.exe",
        "beast uncut leather .mpg.exe",
        "black handjob beast big wifey .avi.exe",
        "black handjob hardcore sleeping feet femdom  (Karin).mpeg.exe",
        "brasilian action beast uncut feet  (Sonja,Tatjana).rar.exe",
        "brasilian horse gay uncut  (Jade).rar.exe",
        "brasilian porn xxx [bangbus] bondage .mpg.exe",
        "bukkake lesbian  (Curtney).mpeg.exe",
        "danish cum sperm licking .mpeg.exe",
        "danish nude trambling voyeur .avi.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking several models cock .rar.exe",
        "gay [milf] glans .mpg.exe",
        "horse hidden feet bondage .rar.exe",
        "italian gang bang sperm [milf] gorgeoushorny .mpg.exe",
        "japanese cum hardcore [free] .rar.exe",
        "japanese fetish hardcore voyeur shoes  (Christine,Liz).mpeg.exe",
        "lesbian [milf] ash .zip.exe",
        "lingerie licking pregnant .avi.exe",
        "lingerie uncut stockings .avi.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "russian cum beast lesbian pregnant .rar.exe",
        "russian kicking fucking big leather .avi.exe",
        "russian kicking xxx full movie beautyfull .mpg.exe",
        "start.exe",
        "trambling lesbian titts wifey .avi.exe",
        "tyrkish action xxx licking redhair .zip.exe",
        "tyrkish cum horse girls feet stockings .mpeg.exe",
        "tyrkish porn horse uncut cock .zip.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}