# dataset/preprocessed/0a8e6bad14cf21f62225b70c3b95c6dd632b0ff16a37579a9fbb3813d0c3550e.json
{
  "hash": "0a8e6bad14cf21f62225b70c3b95c6dd632b0ff16a37579a9fbb3813d0c3550e",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  italian cumshot blowjob hidden blondie .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\italian cumshot blowjob hidden blondie .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\italian cumshot blowjob hidden blondie .mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\0a8e6bad14cf21f62225b70c3b95c6dd632b0ff16a37579a9fbb3813d0c3550e(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\italian cumshot blowjob hidden blondie .mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  brasilian kicking blowjob licking titts upskirt  (Sarah).avi.exe  [sha256: d1c00dbe21c76c10...]\n  brasilian porn fucking hot (!) feet .mpeg.exe  [sha256: 928924f663a83033...]\n  bukkake masturbation  (Janette).zip.exe  [sha256: 78b3a4b3d37daa63...]\n  bukkake public cock sweet .avi.exe  [sha256: 7a07b68d322163d8...]\n  danish cum fucking catfight  (Tatjana).avi.exe  [sha256: cd13e4dcf476b182...]\n  danish cumshot blowjob lesbian penetration .mpg.exe  [sha256: 2ec841218a7a1afc...]\n  danish fetish bukkake catfight .avi.exe  [sha256: 8759e7b3a0aac11a...]\n  danish handjob beast licking hairy .mpg.exe  [sha256: 836e3195c681667a...]\n  danish kicking bukkake hot (!) .rar.exe  [sha256: 4c116f00a5951008...]\n  danish porn beast [milf] granny .avi.exe  [sha256: 992865fe1c600f2e...]\n  danish porn lesbian [bangbus] titts young  (Sarah).mpeg.exe  [sha256: 3de09a71b58cebf6...]\n  defrag64.exe  [sha256: 1c934236d737f517...]\n  deinstall.exe  [sha256: 8f520aee9b0aa238...]\n  drwatson32.exe  [sha256: 47de5bb1cc71043d...]\n  fucking hidden glans .avi.exe  [sha256: f7a2ec9d7688c63e...]\n  gay sleeping hole redhair  (Janette).zip.exe  [sha256: 3b0f803003df063a...]\n  horse big  (Samantha).rar.exe  [sha256: 0fd2ef94c75863c0...]\n  horse big circumcision .zip.exe  [sha256: 95520bcb8c8522d4...]\n  horse uncut titts beautyfull .avi.exe  [sha256: 8840e110a4011350...]\n  italian gang bang hardcore public cock .zip.exe  [sha256: 4c44c1e074854687...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\danish porn lesbian [bangbus] titts young  (Sarah).mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\italian nude lingerie sleeping .rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\horse big circumcision .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\russian porn horse sleeping pregnant .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\xxx public .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\horse big  (Samantha).rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\japanese kicking gay lesbian ash .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\tyrkish kicking gay lesbian  (Liz).mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\italian gang bang hardcore public cock .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\danish cumshot blowjob lesbian penetration .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\danish porn beast [milf] granny .avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\japanese cum bukkake [milf]  (Curtney).zip.exe\n  C:\\WINDOWS\\Temp\\horse uncut titts beautyfull .avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\danish cum fucking catfight  (Tatjana).avi.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\danish kicking bukkake hot (!) .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\gay sleeping hole redhair  (Janette).zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\bukkake masturbation  (Janette).zip.exe\n  C:\\WINDOWS\\ime\\shared\\bukkake public cock sweet .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\japanese cum horse several models lady  (Kathrin,Janette).rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.880077362060547  size=664930\n",
  "ground_truth": {
    "hash": "0a8e6bad14cf21f62225b70c3b95c6dd632b0ff16a37579a9fbb3813d0c3550e",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\trambling beast licking .mpeg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\sperm [free] circumcision .mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\animal big .avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\danish blowjob sleeping  (Sarah).avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\japanese cumshot [free] granny .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\fucking full movie .rar.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\swedish hardcore gang bang catfight .mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\german fetish gang bang several models swallow  (Janette).mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\malaysia fetish handjob girls 40+ .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\sperm girls legs black hairunshaved .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "brasilian kicking blowjob licking titts upskirt  (Sarah).avi.exe",
        "brasilian porn fucking hot (!) feet .mpeg.exe",
        "bukkake masturbation  (Janette).zip.exe",
        "bukkake public cock sweet .avi.exe",
        "danish cum fucking catfight  (Tatjana).avi.exe",
        "danish cumshot blowjob lesbian penetration .mpg.exe",
        "danish fetish bukkake catfight .avi.exe",
        "danish handjob beast licking hairy .mpg.exe",
        "danish kicking bukkake hot (!) .rar.exe",
        "danish porn beast [milf] granny .avi.exe",
        "danish porn lesbian [bangbus] titts young  (Sarah).mpeg.exe",
        "defrag64.exe",
        "deinstall.exe",
        "drwatson32.exe",
        "fucking hidden glans .avi.exe",
        "gay sleeping hole redhair  (Janette).zip.exe",
        "horse big  (Samantha).rar.exe",
        "horse big circumcision .zip.exe",
        "horse uncut titts beautyfull .avi.exe",
        "italian gang bang hardcore public cock .zip.exe",
        "italian nude lingerie sleeping .rar.exe",
        "japanese cum bukkake [milf]  (Curtney).zip.exe",
        "japanese cum horse several models lady  (Kathrin,Janette).rar.exe",
        "japanese kicking gay lesbian ash .zip.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "russian handjob hardcore girls hole .rar.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}