# dataset/preprocessed/0ab0dff3b922e6285e09aa53347868dccf8bcbb165ec79f83084ef08793d3487.json
{
  "hash": "0ab0dff3b922e6285e09aa53347868dccf8bcbb165ec79f83084ef08793d3487",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (1):\n  GET http://carder.bit/  UA: Mozilla/5.0 (Windows NT 6.1; WOW64; Trident/7.0; rv:11.0) like Gecko\n\nDNS LOOKUPS (2):\n  ipv4bot.whatismyipaddress.com \u2192 66.171.248.178\n  ns1.wowservers.ru \u2192 185.103.243.59\n\nIP TRAFFIC (2):\n  185.103.243.59:53 (UDP)\n  66.171.248.178:80 (TCP)\n\nSERVICES OPENED (2):\n  RASMAN\n  Sens\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE (flat):\n  nslookup carder.bit ns1.wowservers.ru\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (6):\n  - contain obfuscated stackstrings\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - modify service\n  - open clipboard\n  - get common file path\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - contain loop\n  - Process::Terminate Process [C0018]\n  - Process::Allocate Thread Local Storage [C0040]\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - modify service\n  - open clipboard\n  - Discovery::File and Directory Discovery [E1083]\n  - linked against CPP standard library\n  - contain a resource (.rsrc) section\n\nHIGHLIGHTED CALLS:\n  SetFileTime\n\n--- ARTIFACTS ---\n\nFILES DROPPED (2):\n  tdeieu.exe  [sha256: 3bd87b8f60f6a220...]\n  C:\\Users\\<USER>\\AppData\\Roaming\\Microsoft\\wmqhms.exe  [sha256: 4a4548ba78fbd475...]\n\nFILES WRITTEN (1):\n  C:\\Documents and Settings\\<USER>\\Application Data\\Microsoft\\tdeieu.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\pjkvlofzefj = \"\n\nOTHER REGISTRY KEYS SET (2):\n  HKU\\S-1-5-18\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\ProxyEnable\n  HKU\\S-1-5-18\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\SavedLegacySettings\n\nMUTEXES CREATED (3):\n  Global\\pc_group=WORKGROUP&ransom_id=ea2177a071d2c847\n  RasPbFile\n  ShimCacheMutex\n\nMUTEXES OPENED (4):\n  Global\\pc_group=WORKGROUP&ransom_id=e8b0adc86c78a9c3\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!history!history.ie5!\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  Local\\c:!users!<USER>!appdata!roaming!microsoft!windows!cookies!\n\nPROCESSES CREATED (1):\n  nslookup carder.bit ns1.wowservers.ru\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, LoadLibraryW, VirtualProtect\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n\nPE SECTIONS (entropy):\n  .text        entropy=6.61\n  .rdata       entropy=5.83\n  .data        entropy=4.12\n  .rsrc        entropy=7.8\n  .reloc       entropy=1.3\n  [overlay]    entropy=2.503258466720581  size=9\n",
  "ground_truth": {
    "hash": "0ab0dff3b922e6285e09aa53347868dccf8bcbb165ec79f83084ef08793d3487",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1115",
        "signature_description": "open clipboard"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1543",
        "signature_description": "modify service"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "contain obfuscated stackstrings",
      "contain loop",
      "terminate process",
      "allocate thread local storage",
      "link function at runtime on Windows",
      "link many functions at runtime",
      "modify service",
      "open clipboard",
      "get common file path",
      "linked against CPP standard library",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [
        "185.103.243.59",
        "66.171.248.178"
      ],
      "domains": [
        "ipv4bot.whatismyipaddress.com",
        "ns1.wowservers.ru"
      ],
      "urls": [
        "http://carder.bit/"
      ],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\pjkvlofzefj"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Users\\<USER>\\AppData\\Roaming\\Microsoft\\wmqhms.exe",
        "tdeieu.exe"
      ],
      "mutexes": [
        "Global\\pc_group=WORKGROUP&ransom_id=ea2177a071d2c847",
        "RasPbFile",
        "ShimCacheMutex"
      ],
      "processes": [
        "nslookup carder.bit ns1.wowservers.ru"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}