# dataset/preprocessed/0aded9f2cddf06823c349e45cb006e416b5482cf7ceb677115cd7f92262dec23.json
{
  "hash": "0aded9f2cddf06823c349e45cb006e416b5482cf7ceb677115cd7f92262dec23",
  "report_text": "--- NETWORK ACTIVITY ---\n\nSERVICES OPENED (1):\n  Messenger\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  bukkake girls feet .avi.exe\n    C:\\Users\\<USER>\\Downloads\\bukkake girls feet .avi.exe\n    C:\\Users\\<USER>\\Downloads\\bukkake girls feet .avi.exe\n      C:\\Users\\<USER>\\Downloads\\bukkake girls feet .avi.exe\n\nCOMMAND EXECUTIONS (1):\n  C:\\Users\\<USER>\\Downloads\\bukkake girls feet .avi.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (17):\n  C:\\Windows\\mssrv.exe  [sha256: 2d62612d47c29d83...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\spanish blowjob kicking hot (!) .mpg.exe  [sha256: 190d9eee3856fc38...]\n  C:\\Program Files\\dotnet\\shared\\beastiality several models cock .zip.exe  [sha256: 2040a5f92385dd70...]\n  C:\\Program Files\\DVD Maker\\Shared\\action action uncut feet shower .zip.exe  [sha256: 9babb5f5b996cb11...]\n  C:\\Program Files\\Windows Journal\\Templates\\african hardcore [free] mistress .rar.exe  [sha256: f2c235b92b570a3a...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\russian beast nude lesbian high heels  (Gina).mpg.exe  [sha256: ca577d871fd08adf...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\russian handjob [free] ash  (Britney).avi.exe  [sha256: ed43ab19646727fc...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\british hardcore hidden .rar.exe  [sha256: d894ac1642efd1e5...]\n  C:\\Program Files (x86)\\dotnet\\shared\\french gay gang bang big mature .mpeg.exe  [sha256: 2d5fc0dab8b83bea...]\n  C:\\Program Files (x86)\\Google\\Update\\Download\\japanese gay nude masturbation traffic .mpg.exe  [sha256: 8ab2d2777e902dab...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\blowjob lesbian feet beautyfull .avi.exe  [sha256: 6285c81b5cdfdfc9...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\indian gay big .zip.exe  [sha256: 078f622e323940c0...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\cumshot [bangbus] legs  (Liz,Kathrin).mpeg.exe  [sha256: fbf4f6b8ecb19e67...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\russian gay voyeur legs .mpeg.exe  [sha256: 9f88ebfc4519da75...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\black beast [milf] ash sm  (Melissa).mpg.exe  [sha256: 5b28045e7d60a8c9...]\n  C:\\Program Files (x86)\\Steam\\clientui\\images\\sharedui\\african fetish handjob masturbation glans  (Curtney).zip.exe  [sha256: 7c7e53d3a84274ff...]\n  C:\\Program Files (x86)\\Steam\\controller_base\\templates\\danish kicking fucking voyeur femdom .mpg.exe  [sha256: 67c15c5b57d28fcc...]\n\nFILES WRITTEN (2):\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\american handjob beastiality catfight ash Y\u008dE\u00e2PS\u00e8&\u0007 .rar.exe\n  C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\brasilian gay nude full movie \u00a4\u00e3\u0018  (Britney,Tatjana).avi.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.879044055938721  size=78512\n",
  "ground_truth": {
    "hash": "0aded9f2cddf06823c349e45cb006e416b5482cf7ceb677115cd7f92262dec23",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9437,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\russian handjob [free] ash  (Britney).avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\british hardcore hidden .rar.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\japanese gay nude masturbation traffic .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\blowjob lesbian feet beautyfull .avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\indian gay big .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\russian gay voyeur legs .mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\cumshot [bangbus] legs  (Liz,Kathrin).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\black beast [milf] ash sm  (Melissa).mpg.exe",
        "C:\\Program Files (x86)\\Steam\\clientui\\images\\sharedui\\african fetish handjob masturbation glans  (Curtney).zip.exe",
        "C:\\Program Files (x86)\\Steam\\controller_base\\templates\\danish kicking fucking voyeur femdom .mpg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\french gay gang bang big mature .mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\spanish blowjob kicking hot (!) .mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\action action uncut feet shower .zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\african hardcore [free] mistress .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\russian beast nude lesbian high heels  (Gina).mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\beastiality several models cock .zip.exe",
        "C:\\Windows\\mssrv.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Jujubox"
  ]
}