# dataset/preprocessed/0b148b1b66e71b7ea3cd94b6664900596ea8223f17e417beb693043863203974.json
{
  "hash": "0b148b1b66e71b7ea3cd94b6664900596ea8223f17e417beb693043863203974",
  "report_text": "--- NETWORK ACTIVITY ---\n\nSERVICES OPENED (1):\n  Messenger\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (2):\n  C:\\Users\\<USER>\\Downloads\\0b148b1b66e71b7e_swedish animal bukkake licking glans swallow .zip.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (10):\n  C:\\Windows\\mssrv.exe  [sha256: fe880349d49becb8...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\cum hardcore licking boobs bondage .zip.exe  [sha256: 7398396c73a3059d...]\n  C:\\Program Files\\dotnet\\shared\\russian blowjob lesbian several models latex .mpeg.exe  [sha256: 03135ebaee2477f7...]\n  C:\\Program Files\\DVD Maker\\Shared\\blowjob sleeping redhair .mpeg.exe  [sha256: 168475156b681afc...]\n  C:\\Program Files\\Windows Journal\\Templates\\canadian porn blowjob sleeping redhair .rar.exe  [sha256: 8af92440a5da708f...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\tyrkish beastiality public boobs .mpeg.exe  [sha256: e78b239c3b3fafa0...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\swedish action beast [milf] fishy .rar.exe  [sha256: 2c2781f7a5ca39f7...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\spanish lesbian fucking several models balls  (Sonja,Samantha).mpeg.exe  [sha256: ed09f41236d921a1...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\swedish gang bang big .zip.exe  [sha256: a58d602413933dcf...]\n  C:\\Program Files (x86)\\dotnet\\shared\\asian action catfight vagina lady  (Liz).mpeg.exe  [sha256: 49511166091a5222...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\Administrator\\Application Data\\Macromedia\\Flash Player\\#SharedObjects\\handjob lingerie voyeur beautyfull  (Christine).mpg.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Bin\\bse_temp\\german hardcore fucking [milf] blondie .rar.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Extension\\com.sogou.share\\xxx cum uncut glans .mpeg.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Temp\\japanese cumshot hot (!) penetration .zip.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQPCMgr\\Download\\kicking handjob masturbation .avi.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\Download\\hardcore porn catfight .mpeg.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\horse hidden boots  (Liz,Ashley).mpg.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\blowjob horse lesbian cock wifey .mpg.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\Temp\\fetish handjob hidden femdom  (Sylvia).rar.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\TXSSO\\SSOTemp\\tyrkish beast catfight high heels  (Samantha,Jade).mpg.exe\n  C:\\Documents and Settings\\Administrator\\CMB\\PB40\\Temp\\asian gay licking hotel .mpeg.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\gay voyeur .rar.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temporary Internet Files\\horse [free] .rar.exe\n  C:\\Documents and Settings\\Administrator\\Templates\\black trambling lingerie sleeping vagina beautyfull  (Ashley,Sylvia).zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\sperm full movie beautyfull .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\asian nude uncut .rar.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temporary Internet Files\\danish animal hidden ash .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\italian action several models ejaculation .mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\animal big traffic  (Anniston).avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temporary Internet Files\\gang bang handjob licking .zip.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880218505859375  size=1174665\n",
  "ground_truth": {
    "hash": "0b148b1b66e71b7ea3cd94b6664900596ea8223f17e417beb693043863203974",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9444,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "compiled with Go",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\swedish action beast [milf] fishy .rar.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\spanish lesbian fucking several models balls  (Sonja,Samantha).mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\swedish gang bang big .zip.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\asian action catfight vagina lady  (Liz).mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\cum hardcore licking boobs bondage .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\blowjob sleeping redhair .mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\canadian porn blowjob sleeping redhair .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\tyrkish beastiality public boobs .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\russian blowjob lesbian several models latex .mpeg.exe",
        "C:\\Windows\\mssrv.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Jujubox"
  ]
}