# dataset/preprocessed/0c0dea9e841c8812b8856ebea9f4457e27e76dde6c36632697ebf38d28c51b55.json
{
  "hash": "0c0dea9e841c8812b8856ebea9f4457e27e76dde6c36632697ebf38d28c51b55",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (22):\n  GET http://uwsozi.net/244/510.html\n  GET http://uwsozi.net/224/167.html\n  GET http://mkkuei4kdsz.com/913/981.html\n  GET http://ow5dirasuek.com/152/518.html\n  GET http://uwsozi.net/847/524.html\n  GET http://uwsozi.net/585/184.html\n  GET http://mkkuei4kdsz.com/484/284.html\n  GET http://ow5dirasuek.com/799/990.html\n  GET http://uwsozi.net/997/181.html\n  GET http://uwsozi.net/960/658.html\n  GET http://mkkuei4kdsz.com/33/235.html\n  GET http://ow5dirasuek.com/557/989.html\n  GET http://uwsozi.net/316/494.html\n  GET http://uwsozi.net/113/892.html\n  GET http://mkkuei4kdsz.com/599/106.html\n\nDNS LOOKUPS (3):\n  mkkuei4kdsz.com\n  ow5dirasuek.com\n  uwsozi.net\n\nSERVICES OPENED (2):\n  RASMAN\n  Sens\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  l4r1sx.exe\n    C:\\Users\\<USER>\\AppData\\Roaming\\ioylax84.exe\n\nCOMMAND EXECUTIONS (3):\n  ( 1224 - cmd.exe ) c:\\windows\\system32\\cmd.exe /c \"c:\\users\\user41~1\\appdata\\local\\temp\\gwtrtbs.exe\"\n  ( 1016 - gwtrtbs.exe ) c:\\users\\user41~1\\appdata\\local\\temp\\gwtrtbs.exe\n  ( 1384 - juoyn3g1.exe ) c:\\users\\USER\\appdata\\roaming\\juoyn3g1.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (17):\n  - encode data using XOR\n  - contain obfuscated stackstrings\n  - check Internet connectivity via WinINet\n  - get common file path\n  - link function at runtime on Windows\n  - encrypt data using RC4 PRGA\n  - query environment variable\n  - query or enumerate registry key\n  - query or enumerate registry value\n  - download and write a file\n  - get file size\n  - enumerate files on Windows\n  - enumerate files recursively\n  - persist via AppInit_DLLs registry key\n  - find graphical window\n  - accept command line arguments\n  - persist via Run registry key\n\nCAPA RULE MATCHES (MBC):\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Process [C0017]\n  - Memory::Allocate Memory [C0007]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - File System::Create File [C0016]\n  - check Internet connectivity via WinINet\n  - query remote server for available data\n  - reference absolute stream path on Windows\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - contain loop\n  - File System::Delete File [C0047]\n  - Discovery::File and Directory Discovery [E1083]\n  - link function at runtime on Windows\n  - Cryptography::Encrypt Data::RC4 [C0027.009]\n  - Cryptography::Generate Pseudo-random Sequence::RC4 PRGA [C0021.004]\n  - File System::Read File [C0051]\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Key [C0036.005]\n\nHIGHLIGHTED CALLS:\n  SetFileTime\n\n--- ARTIFACTS ---\n\nFILES DROPPED (3):\n  l4r1sx.exe  [sha256: 19eeb09725f60873...]\n  C:\\Users\\<USER>\\AppData\\Roaming\\ioylax84.exe  [sha256: dd955c9ade61b22f...]\n  c:\\users\\USER\\appdata\\roaming\\juoyn3g1.exe  [sha256: 7b173891081344a2...]\n\nFILES WRITTEN (5):\n  ( 1016 - gwtrtbs.exe ) c:\\users\\USER\\appdata\\roaming\\juoyn3g1.exe\n  C:\\Documents and Settings\\<USER>\\Application Data42wwmnt8u\n  C:\\Documents and Settings\\<USER>\\Application Data\\l4r1sx.exe\n  C:\\Documents and Settings\\Administrator\\Application Data0v3k\n  C:\\Documents and Settings\\Administrator\\Application Data\\7t5abxq.exe\n\nPERSISTENCE REGISTRY KEYS (3):\n  HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\0jzawm = C\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\8ajzv5 = C:\\Documents and Settings\\Administrator\\Application Data\\7t5abxq.exe\u0000\n  \\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\run\\ikgy = c:\\users\\USER\\appdata\\roaming\\juoyn3g1.exe\n\nOTHER REGISTRY KEYS SET (9):\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\ProxyEnable\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\n  \\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\internet setting\n  \\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\internet setting\n  \\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\internet setting\n  \\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\internet setting\n  \\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\internet setting\n\nMUTEXES CREATED (13):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  IESQMMUTEX_0_208\n  RasPbFile\n  Semaphore_157862\n  ShimCacheMutex\n\nMUTEXES OPENED (7):\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!history!history.ie5!\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  Local\\c:!users!<USER>!appdata!roaming!microsoft!windows!cookies!\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\<USER>\\Application Data\\l4r1sx.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessW, GetProcAddress, LoadLibraryW, VirtualAlloc, VirtualProtect\n\nCOMPILER AND PACKER:\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n\nPE SECTIONS (entropy):\n  .text        entropy=6.64\n  .rdata       entropy=4.26\n  .data        entropy=3.66\n  [overlay]    entropy=1.4528098106384277  size=248\n",
  "ground_truth": {
    "hash": "0c0dea9e841c8812b8856ebea9f4457e27e76dde6c36632697ebf38d28c51b55",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8889,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry key"
      },
      {
        "id": "T1016",
        "signature_description": "check Internet connectivity via WinINet"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "query environment variable"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1105",
        "signature_description": "download and write a file"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1546",
        "signature_description": "persist via AppInit_DLLs registry key"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      }
    ],
    "capa_sig_names": [
      "delay execution",
      "create process on Windows",
      "allocate memory",
      "generate random numbers using the Delphi LCG",
      "encode data using XOR",
      "create or open registry key",
      "contain obfuscated stackstrings",
      "create or open file",
      "check Internet connectivity via WinINet",
      "query remote server for available data",
      "allocate RW memory",
      "reference absolute stream path on Windows",
      "create thread",
      "check mutex",
      "contain loop",
      "delete file",
      "get common file path",
      "link function at runtime on Windows",
      "encrypt data using RC4 PRGA",
      "read .ini file",
      "query environment variable",
      "query or enumerate registry key",
      "query or enumerate registry value",
      "set registry value",
      "get OS version",
      "create HTTP request",
      "connect to URL",
      "read data from Internet",
      "receive data",
      "write file on Windows",
      "download and write a file",
      "receive and write data from server to client",
      "get file size",
      "read file on Windows",
      "hash data with MD5",
      "enumerate files on Windows",
      "enumerate files recursively",
      "persist via AppInit_DLLs registry key",
      "send HTTP request",
      "connect to HTTP server",
      "send data",
      "check for time delay via GetTickCount",
      "find graphical window",
      "find taskbar",
      "accept command line arguments",
      "persist via Run registry key",
      "terminate process",
      "create mutex",
      "copy file",
      "check mutex and exit"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "mkkuei4kdsz.com",
        "ow5dirasuek.com",
        "uwsozi.net"
      ],
      "urls": [
        "http://mkkuei4kdsz.com/114/285.html",
        "http://mkkuei4kdsz.com/33/235.html",
        "http://mkkuei4kdsz.com/484/284.html",
        "http://mkkuei4kdsz.com/599/106.html",
        "http://mkkuei4kdsz.com/908/152.html",
        "http://mkkuei4kdsz.com/913/981.html",
        "http://ow5dirasuek.com/152/518.html",
        "http://ow5dirasuek.com/222/50.html",
        "http://ow5dirasuek.com/557/989.html",
        "http://ow5dirasuek.com/735/452.html",
        "http://ow5dirasuek.com/799/990.html",
        "http://uwsozi.net/113/892.html",
        "http://uwsozi.net/224/167.html",
        "http://uwsozi.net/244/510.html",
        "http://uwsozi.net/316/494.html",
        "http://uwsozi.net/425/947.html",
        "http://uwsozi.net/585/184.html",
        "http://uwsozi.net/600/636.html",
        "http://uwsozi.net/847/524.html",
        "http://uwsozi.net/862/123.html",
        "http://uwsozi.net/960/658.html",
        "http://uwsozi.net/997/181.html"
      ],
      "registry_persistence": [
        "HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\0jzawm",
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\8ajzv5",
        "\\registry\\user\\s-1-5-21-2589874625-3071026754-1203412979-1000\\software\\microsoft\\windows\\currentversion\\run\\ikgy"
      ],
      "registry_all_count": 12,
      "files_dropped": [
        "C:\\Users\\<USER>\\AppData\\Roaming\\ioylax84.exe",
        "c:\\users\\USER\\appdata\\roaming\\juoyn3g1.exe",
        "l4r1sx.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "IESQMMUTEX_0_208",
        "RasPbFile",
        "Semaphore_157862",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "C:\\Documents and Settings\\<USER>\\Application Data\\l4r1sx.exe"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "ReaQta-Hive",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}