# dataset/preprocessed/0c262dce428e928a567c6bef7cfea3328d0392b2990b9b84226742998a5ce38e.json
{
  "hash": "0c262dce428e928a567c6bef7cfea3328d0392b2990b9b84226742998a5ce38e",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (3):\n  ((null)) C:\\0c262dce428e928a567c6bef7cfea3328d0392b2990b9b84226742998a5ce38e(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\horse [milf] .avi.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african horse big .zip.exe  [sha256: ff5d0ada17768397...]\n  african nude lesbian [free] sm .mpg.exe  [sha256: 7cfacb24d5b63b1d...]\n  american porn [bangbus] glans sweet .zip.exe  [sha256: 559dfa10f926ce35...]\n  animal sleeping latex .zip.exe  [sha256: 9563df5ef8940936...]\n  asian action public circumcision  (Britney,Janette).mpg.exe  [sha256: f044fa45ab0ee56b...]\n  beastiality beastiality [bangbus] legs stockings .mpg.exe  [sha256: 1d25de71f11a0ecc...]\n  black xxx blowjob girls legs black hairunshaved .rar.exe  [sha256: 75cb9f5f659f3588...]\n  blowjob hardcore licking blondie .mpeg.exe  [sha256: fce27e3494c26041...]\n  bukkake uncut leather .avi.exe  [sha256: 27ddad729c153ae4...]\n  canadian trambling [bangbus] feet .mpg.exe  [sha256: fab3893fe7de737e...]\n  cum gay [milf] .zip.exe  [sha256: 533e1cbeee2d469b...]\n  cum masturbation 50+ .rar.exe  [sha256: b271b125e120acf5...]\n  danish fetish sleeping hotel  (Jade,Sonja).rar.exe  [sha256: 2f146c8cf3871051...]\n  defrag.exe  [sha256: 5d4673c92547f46b...]\n  driver.exe  [sha256: 23f4f2467d0134a9...]\n  edit.exe  [sha256: eb8297bbc63d19a8...]\n  freeftp.exe  [sha256: 7ea408c8ecb03720...]\n  french beast licking  (Sandy).mpg.exe  [sha256: 6934cd5c6ecdcc1f...]\n  french horse sperm licking .zip.exe  [sha256: 7324611d884cfc0b...]\n  french sperm masturbation  (Anniston,Ashley).rar.exe  [sha256: 56f7b2fd8d075b9d...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\canadian trambling [bangbus] feet .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\cum gay [milf] .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\african nude lesbian [free] sm .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\bukkake uncut leather .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\asian action public circumcision  (Britney,Janette).mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\blowjob hardcore licking blondie .mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\animal sleeping latex .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\italian gay big hole .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\african horse big .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\norwegian lesbian trambling catfight hotel .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\indian sperm uncut legs bondage  (Liz,Curtney).avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\norwegian gay hardcore [free] .rar.exe\n  C:\\WINDOWS\\Temp\\black xxx blowjob girls legs black hairunshaved .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\french horse sperm licking .zip.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\indian xxx action [free] boots .mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\trambling lingerie hidden mistress .zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\french beast licking  (Sandy).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\gay big titts .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\cum masturbation 50+ .rar.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.875197410583496  size=12651\n",
  "ground_truth": {
    "hash": "0c262dce428e928a567c6bef7cfea3328d0392b2990b9b84226742998a5ce38e",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\brasilian beastiality hardcore [bangbus]  (Sarah).mpeg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\xxx several models  (Jade).mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\bukkake sleeping .rar.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\sperm licking latex .mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\brasilian cum lingerie catfight feet fishy .mpg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\trambling licking feet stockings .rar.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\bukkake masturbation titts redhair .avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\hardcore big hole traffic  (Tatjana).avi.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\swedish gang bang xxx licking hole .mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\danish cum blowjob uncut feet mature  (Tatjana).avi.exe",
        "C:\\Program Files\\dotnet\\shared\\beast sleeping shoes .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "african horse big .zip.exe",
        "african nude lesbian [free] sm .mpg.exe",
        "american porn [bangbus] glans sweet .zip.exe",
        "animal sleeping latex .zip.exe",
        "asian action public circumcision  (Britney,Janette).mpg.exe",
        "beastiality beastiality [bangbus] legs stockings .mpg.exe",
        "black xxx blowjob girls legs black hairunshaved .rar.exe",
        "blowjob hardcore licking blondie .mpeg.exe",
        "bukkake uncut leather .avi.exe",
        "canadian trambling [bangbus] feet .mpg.exe",
        "cum gay [milf] .zip.exe",
        "cum masturbation 50+ .rar.exe",
        "danish fetish sleeping hotel  (Jade,Sonja).rar.exe",
        "defrag.exe",
        "driver.exe",
        "edit.exe",
        "freeftp.exe",
        "french beast licking  (Sandy).mpg.exe",
        "french horse sperm licking .zip.exe",
        "french sperm masturbation  (Anniston,Ashley).rar.exe",
        "fucking horse [bangbus] hole balls  (Janette).rar.exe",
        "game3d.exe",
        "gay big titts .rar.exe",
        "hdcleaner.exe",
        "indian nude several models pregnant .zip.exe",
        "indian sperm uncut legs bondage  (Liz,Curtney).avi.exe",
        "indian xxx action [free] boots .mpg.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}