# dataset/preprocessed/0f85ab76271d4ecfebe8414b0c27a18f61f9c8e418c5cfad35573a7d89dab1c8.json
{
  "hash": "0f85ab76271d4ecfebe8414b0c27a18f61f9c8e418c5cfad35573a7d89dab1c8",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (5):\n  ftp.0adz.com \u2192 45.33.23.183\n  ftp.0adz.com \u2192 45.79.19.196, 198.58.118.167\n  ftp.search.instantfreesite.com\n  gripto.leadhoster.com\n  gripto.leadhoster.com \u2192 127.0.0.1\n\nIP TRAFFIC (3):\n  198.58.118.167:21 (TCP)\n  45.79.19.196:21 (TCP)\n  8.8.8.8:53 (UDP)\n\nSERVICES OPENED (1):\n  RASMAN\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  jusched.exe\n    C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n\nCOMMAND EXECUTIONS (8):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n  \"C:\\Program Files\\d41b8753\\jusched.exe\"\n  C:\\Windows\\system32\\schtasks.exe /delete /f /TN \"Microsoft\\Windows\\Customer Experience Improvement Program\\Uploader\"\n  C:\\Windows\\System32\\wsqmcons.exe \n  \\??\\C:\\Windows\\system32\\conhost.exe\n  \"C:\\Program Files\\48db4bad\\jusched.exe\" \n   \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - enumerate files on Windows\n  - parse PE header\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - Discovery::File and Directory Discovery [E1083]\n  - parse PE header\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - contain a resource (.rsrc) section\n  - debug build\n\n--- ARTIFACTS ---\n\nFILES DROPPED (3):\n  jusched.exe  [sha256: 283a95e6fe92b13a...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: 4a480d7a9afe66a2...]\n  \\Program Files\\d41b8753\\jusched.exe  [sha256: fdd7de922ca67d98...]\n\nFILES WRITTEN (17):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files (x86)\\71d2c87e\\info_a\n  C:\\Program Files\\48db4bad\n  C:\\Program Files\\48db4bad\\48db4bad\n  C:\\Program Files\\48db4bad\\info_a\n  C:\\Program Files\\48db4bad\\jusched.exe\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\info_a\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\Program Files\\d41b8753\\d41b8753\n  C:\\Program Files\\d41b8753\\info_a\n  C:\\Program Files\\d41b8753\\jusched.exe\n  C:\\Program Files\\f440491e\\f440491e\n  C:\\Program Files\\f440491e\\info_a\n  C:\\Program Files\\f440491e\\jusched.exe\n  C:\\WINDOWS\\Tasks\\Update23.job\n  C:\\Windows\\Tasks\\Update23.job\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:JavaUpdate23\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe = C:\\Program Files\\f440491e\\jusched.exe:*:Enabled:JavaUpdate23\u0000\n\nOTHER REGISTRY KEYS SET (10):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Program \n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\n  \\REGISTRY\\USER\\S-1-5-21-3616907901-860437140-2684471051-1001\\HexStream = 53006F006600740077006100720065005C004D006900630\n  \\REGISTRY\\USER\\S-1-5-21-3616907901-860437140-2684471051-1001\\HexStream = 53006F006600740077006100720065005C004D006900630\n  \\REGISTRY\\USER\\S-1-5-21-3616907901-860437140-2684471051-1001\\HexStream = 53006F006600740077006100720065005C004D006900630\n  HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\SQMClient\\Windows\\AdaptiveSqm\\ManifestInfo\\Version\n  HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\SQMClient\\Windows\\WSqmConsLastEventTimeStamp\n  HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\SQMClient\\Windows\\WSqmConsLastRunTime\n  HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\AutoDetect\n  HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\UNCAsIntranet\n\nMUTEXES CREATED (24):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  Global\\SQMWindowsConsolidator\n  Local\\MSCTF.Asm.MutexDefault0\n  Local\\WininetConnectionMutex\n  Local\\WininetProxyRegistryMutex\n\nMUTEXES OPENED (4):\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\nPROCESSES CREATED (4):\n  C:\\Program Files\\48db4bad\\jusched.exe\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n  C:\\Windows\\System32\\schtasks.exe\n  jusched.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.92\n  .rdata       entropy=4.65\n  .data        entropy=0.87\n  .idata       entropy=3.01\n  .rsrc        entropy=3.08\n  .reloc       entropy=4.25\n  [overlay]    entropy=6.585947513580322  size=6739\n",
  "ground_truth": {
    "hash": "0f85ab76271d4ecfebe8414b0c27a18f61f9c8e418c5cfad35573a7d89dab1c8",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9167,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "enumerate files on Windows",
      "parse PE header",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "contain a resource (.rsrc) section",
      "debug build"
    ],
    "iocs": {
      "ips": [
        "198.58.118.167",
        "45.79.19.196",
        "8.8.8.8"
      ],
      "domains": [
        "ftp.0adz.com",
        "ftp.search.instantfreesite.com",
        "gripto.leadhoster.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe"
      ],
      "registry_all_count": 12,
      "files_dropped": [
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "\\Program Files\\d41b8753\\jusched.exe",
        "jusched.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "Global\\SQMWindowsConsolidator",
        "Local\\MSCTF.Asm.MutexDefault0",
        "Local\\WininetConnectionMutex",
        "Local\\WininetProxyRegistryMutex",
        "Local\\WininetStartupMutex",
        "Local\\ZoneAttributeCacheCounterMutex",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex",
        "Local\\_!MSFTHISTORY!_",
        "Local\\c:!users!<USER>!appdata!local!microsoft!windows!history!history.ie5!",
        "Local\\c:!users!<USER>!appdata!local!microsoft!windows!temporary internet files!content.ie5!",
        "Local\\c:!users!<USER>!appdata!roaming!microsoft!windows!cookies!",
        "RasPbFile",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "C:\\Program Files\\48db4bad\\jusched.exe",
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"",
        "C:\\Windows\\System32\\schtasks.exe",
        "jusched.exe"
      ]
    }
  },
  "num_sandboxes": 6,
  "sandbox_names": [
    "CAPA",
    "Rising MOVES",
    "SNDBOX",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}