# dataset/preprocessed/0fc44670392d0bbf00e50136fd6b90f1c74e8b6754581edc7e9edd1aa6a4d029.json
{
  "hash": "0fc44670392d0bbf00e50136fd6b90f1c74e8b6754581edc7e9edd1aa6a4d029",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  0fc44670392d0bbf_mssrv.exe\n    C:\\Users\\<USER>\\Downloads\\0fc44670392d0bbf_mssrv.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\0fc44670392d0bbf00e50136fd6b90f1c74e8b6754581edc7e9edd1aa6a4d029(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\0fc44670392d0bbf_mssrv.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  action masturbation .rar.exe  [sha256: 59e3ab4b4075b03c...]\n  african trambling [bangbus] .avi.exe  [sha256: ede3d39b3abe33e8...]\n  american fetish beastiality girls hole  (Christine).avi.exe  [sha256: d7c748e373d6ee86...]\n  american horse cum sleeping granny .mpeg.exe  [sha256: 17ed5f4666ec5310...]\n  blowjob animal [milf]  (Kathrin).avi.exe  [sha256: 4e6224b5f6a42716...]\n  blowjob trambling sleeping sweet .mpeg.exe  [sha256: 046b9f58a136c55f...]\n  chinese cumshot animal hot (!) .avi.exe  [sha256: fe61cf4f7316d35a...]\n  cumshot action public mature .mpeg.exe  [sha256: af992ee432d5984b...]\n  danish bukkake catfight .rar.exe  [sha256: 15cf939792ec897a...]\n  danish fucking lingerie masturbation boobs circumcision  (Gina).zip.exe  [sha256: 0d889e61e6269213...]\n  driver.exe  [sha256: a6759c6374e6d81b...]\n  fetish porn several models pregnant  (Christine).mpeg.exe  [sha256: a776e7204a128c6d...]\n  french sperm [bangbus] mature  (Kathrin).mpg.exe  [sha256: 680decaf1b605700...]\n  fucking blowjob [milf] .mpg.exe  [sha256: b09e708d79bf79b3...]\n  german beastiality gang bang uncut balls  (Britney,Christine).mpg.exe  [sha256: 29d36ef814a88398...]\n  german bukkake beast uncut cock .mpeg.exe  [sha256: 47ff4ee28702b2bf...]\n  horse bukkake [free] penetration .avi.exe  [sha256: 69ed0fd94326d3d2...]\n  horse hot (!) girly  (Britney).rar.exe  [sha256: 3c1489fc19368831...]\n  italian horse catfight .mpg.exe  [sha256: a001f4439b5405cc...]\n  japanese gang bang sleeping girly .avi.exe  [sha256: 36297c1caee678e9...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\russian blowjob nude girls gorgeoushorny .rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\japanese gang bang sleeping girly .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\russian beastiality licking ash balls .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\american horse cum sleeping granny .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\xxx [free] nipples .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\danish bukkake catfight .rar.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\cumshot action public mature .mpeg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\blowjob trambling sleeping sweet .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\american fetish beastiality girls hole  (Christine).avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\chinese cumshot animal hot (!) .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\kicking [bangbus] swallow  (Ashley,Ashley).mpg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\horse hot (!) girly  (Britney).rar.exe\n  C:\\WINDOWS\\Temp\\german bukkake beast uncut cock .mpeg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\african trambling [bangbus] .avi.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\german beastiality gang bang uncut balls  (Britney,Christine).mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\blowjob animal [milf]  (Kathrin).avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\xxx trambling catfight .mpg.exe\n  C:\\WINDOWS\\ime\\shared\\spanish nude animal hidden balls  (Anniston).rar.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\norwegian fetish kicking [milf] titts .mpg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.87510347366333  size=37266\n",
  "ground_truth": {
    "hash": "0fc44670392d0bbf00e50136fd6b90f1c74e8b6754581edc7e9edd1aa6a4d029",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9437,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\asian trambling fucking catfight gorgeoushorny  (Karin).avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\british bukkake lingerie hot (!) glans granny .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\german horse horse uncut girly  (Christine,Anniston).zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\fetish hidden castration  (Jade).zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\russian bukkake [milf] hole fishy .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\norwegian kicking several models feet .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "action masturbation .rar.exe",
        "african trambling [bangbus] .avi.exe",
        "american fetish beastiality girls hole  (Christine).avi.exe",
        "american horse cum sleeping granny .mpeg.exe",
        "blowjob animal [milf]  (Kathrin).avi.exe",
        "blowjob trambling sleeping sweet .mpeg.exe",
        "chinese cumshot animal hot (!) .avi.exe",
        "cumshot action public mature .mpeg.exe",
        "danish bukkake catfight .rar.exe",
        "danish fucking lingerie masturbation boobs circumcision  (Gina).zip.exe",
        "driver.exe",
        "fetish porn several models pregnant  (Christine).mpeg.exe",
        "french sperm [bangbus] mature  (Kathrin).mpg.exe",
        "fucking blowjob [milf] .mpg.exe",
        "german beastiality gang bang uncut balls  (Britney,Christine).mpg.exe",
        "german bukkake beast uncut cock .mpeg.exe",
        "horse bukkake [free] penetration .avi.exe",
        "horse hot (!) girly  (Britney).rar.exe",
        "italian horse catfight .mpg.exe",
        "japanese gang bang sleeping girly .avi.exe",
        "kicking [bangbus] swallow  (Ashley,Ashley).mpg.exe",
        "mssrv.exe",
        "norwegian fetish kicking [milf] titts .mpg.exe",
        "pgp9.exe",
        "russian beastiality licking ash balls .zip.exe",
        "russian blowjob nude girls gorgeoushorny .rar.exe",
        "russian gang bang sleeping .rar.exe",
        "safetyserver.exe",
        "screensaver.scr",
        "setup.exe",
        "spanish nude animal hidden balls  (Anniston).rar.exe",
        "start.exe",
        "tyrkish fetish masturbation .mpeg.exe",
        "undelete.exe",
        "uninstall.exe",
        "word.exe",
        "xxx [free] nipples .zip.exe",
        "xxx trambling catfight .mpg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}