# dataset/preprocessed/10985933f28a320cbe0696e6bd4414735b46cb01777778e4307c83077d6f06f8.json
{
  "hash": "10985933f28a320cbe0696e6bd4414735b46cb01777778e4307c83077d6f06f8",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  9a35015ef1af2f900dd37136159e6eef.virus.exe\n    C:\\Users\\<USER>\\Downloads\\9a35015ef1af2f900dd37136159e6eef.virus.mm /zhj\n      C:\\Windows\\YZH.exe /zhj\n    C:\\Users\\<USER>\\Downloads\\9a35015ef1af2f900dd37136159e6eef.virus.tmp\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (26):\n  - set file attributes\n  - acquire debug privileges\n  - encode data using XOR\n  - hide graphical window\n  - accept command line arguments\n  - query or enumerate registry value\n  - link function at runtime on Windows\n  - get geographical location\n  - get disk size\n  - link many functions at runtime\n  - get system information on Windows\n  - capture screenshot\n  - check if file exists\n  - find graphical window\n  - modify access privileges\n  - get hostname\n  - get common file path\n  - delete registry value\n  - delete registry key\n  - query or enumerate registry key\n  - encode data using Base64\n  - log keystrokes via application hook\n  - persist via Run registry key\n  - packed with UPX\n  - reference Base64 string\n  - compiled with AutoIt\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - File System::Create File [C0016]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Get File Attributes [C0049]\n  - Process::Open Process [C0065]\n  - acquire debug privileges\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - Process::Create Process::Create Suspended Process [C0017.003]\n  - Process::Create Process [C0017]\n  - hide graphical window\n  - contain loop\n  - File System::Create Directory [C0046]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - File System::Writes File [C0052]\n  - link function at runtime on Windows\n  - get geographical location\n  - Process::Set Thread Local Storage Value [C0041]\n  - get thread local storage value\n\n--- ARTIFACTS ---\n\nFILES DROPPED (1):\n  C:\\Users\\<USER>\\Downloads\\9a35015ef1af2f900dd37136159e6eef.virus.mm  [sha256: f10ca14e0f50249c...]\n\nFILES WRITTEN (7):\n  C:\\Program Files\\Sublime Text 3\\crash_reporter.exe\n  C:\\Program Files\\Sublime Text 3\\plugin_host.exe\n  C:\\Program Files\\Sublime Text 3\\subl.exe\n  C:\\Program Files\\dotnet\\dotnet.exe\n  C:\\Users\\<USER>\\Downloads\\9a35015ef1af2f900dd37136159e6eef.virus.tmp\n  C:\\dev_man.exe\n  C:\\sdel.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SoftWare\\Microsoft\\Windows\\CurrentVersion\\Run\\YZH = C:\\Windows\\YZH.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  advapi32.dll: AdjustTokenPrivileges, OpenProcessToken\n  KERNEL32.DLL: CreateProcessA, GetProcAddress, VirtualAlloc, WinExec\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [sfx] temporary EXE SFX\n  [Packer] UPX\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  [Overlay] EXE file\n\nPE SECTIONS (entropy):\n  UPX0         entropy=6.52\n  UPX1         entropy=3.29\n  .rsrc        entropy=3.11\n  stub         entropy=4.79\n  [overlay]    entropy=6.293686389923096  size=425403\n",
  "ground_truth": {
    "hash": "10985933f28a320cbe0696e6bd4414735b46cb01777778e4307c83077d6f06f8",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8732,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via application hook"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "check if file exists"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry value"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1134",
        "signature_description": "acquire debug privileges"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "create or open registry key",
      "create or open file",
      "set file attributes",
      "delay execution",
      "get file attributes",
      "open process",
      "acquire debug privileges",
      "encode data using XOR",
      "PEB access",
      "create process suspended",
      "create process on Windows",
      "hide graphical window",
      "contain loop",
      "create directory",
      "accept command line arguments",
      "query or enumerate registry value",
      "write file on Windows",
      "link function at runtime on Windows",
      "get geographical location",
      "set thread local storage value",
      "get thread local storage value",
      "delete file",
      "get disk size",
      "get OS version",
      "link many functions at runtime",
      "get system information on Windows",
      "capture screenshot",
      "check if file exists",
      "find graphical window",
      "get graphical window text",
      "terminate process",
      "modify access privileges",
      "get hostname",
      "get common file path",
      "set registry value",
      "delete registry value",
      "delete registry key",
      "query or enumerate registry key",
      "resolve DNS",
      "send data on socket",
      "send data",
      "receive data on socket",
      "receive data",
      "encode data using Base64",
      "copy file",
      "set application hook",
      "log keystrokes via application hook",
      "persist via Run registry key",
      "packed with UPX",
      "contain an embedded PE file",
      "contain a resource (.rsrc) section",
      "reference Base64 string",
      "compiled with AutoIt",
      "compiled with Borland Delphi",
      "(internal) autoit file limitation",
      "(internal) packer file limitation",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SoftWare\\Microsoft\\Windows\\CurrentVersion\\Run\\YZH"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Users\\<USER>\\Downloads\\9a35015ef1af2f900dd37136159e6eef.virus.mm"
      ],
      "mutexes": [],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Jujubox"
  ]
}