# dataset/preprocessed/118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea.json
{
  "hash": "118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  winhash_up.exe\n    C:\\Windows\\bugMAKER.bat\n\nCOMMAND EXECUTIONS (2):\n  C:\\Users\\<USER>\\Downloads\\winhash_up.exe.exe\n  cmd /c C:\\WINDOWS\\bugMAKER.bat\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (8):\n  - set file attributes\n  - accept command line arguments\n  - query or enumerate registry value\n  - link function at runtime on Windows\n  - get geographical location\n  - get disk size\n  - get common file path\n  - persist via Run registry key\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - File System::Create File [C0016]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - reference absolute stream path on Windows\n  - File System::Set File Attributes [C0050]\n  - Process::Create Process [C0017]\n  - contain loop\n  - File System::Create Directory [C0046]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - link function at runtime on Windows\n  - get geographical location\n  - Process::Set Thread Local Storage Value [C0041]\n  - get thread local storage value\n  - File System::Delete File [C0047]\n  - Discovery::System Information Discovery [E1082]\n  - get OS version\n  - Operating System::Registry::Set Registry Key [C0036.001]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (2):\n  winhash_up.exe  [sha256: 1b9f9d2a681865a2...]\n  C:\\Windows\\winhash_up.exe  [sha256: de25a4464ccdec50...]\n\nFILES WRITTEN (27):\n  C:\\BalROK.HTML\n  C:\\Users\\<USER>\\Downloads\\winhash_up.exe.exe\n  C:\\WINDOWS\\SHARE_TEMP\\Icon1.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon10.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon11.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon12.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon13.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon14.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon2.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon5.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon6.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon7.ico\n  C:\\WINDOWS\\SHARE_TEMP\\Icon8.ico\n  C:\\WINDOWS\\bugMAKER.bat\n  C:\\WINDOWS\\winhash_up.exe\n  C:\\WINDOWS\\winhash_up.exez\n  C:\\Windows\\SHARE_TEMP\\Icon10.ico\n  C:\\Windows\\SHARE_TEMP\\Icon12.ico\n  C:\\Windows\\SHARE_TEMP\\Icon13.ico\n  C:\\Windows\\SHARE_TEMP\\Icon14.ico\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft I Service = C:\\Windows\\winhash_up.exe /REGstart\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft I Service = C:\\WINDOWS\\winhash_up.exe /REGstart\u0000\n\nMUTEXES CREATED (1):\n  ShimCacheMutex\n\nPROCESSES CREATED (3):\n  C:\\118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea.exe C:\\118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea.exe\n  C:\\WINDOWS\\bugMAKER.bat\n  C:\\WINDOWS\\system32\\cmd.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  kernel32.dll: GetProcAddress, VirtualAlloc, WinExec\n\nCOMPILER AND PACKER:\n  [sfx] temporary EXE SFX\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  [Overlay] EXE file\n\nPE SECTIONS (entropy):\n  CODE         entropy=6.57\n  DATA         entropy=3.99\n  BSS          entropy=0.0\n  .idata       entropy=4.72\n  .tls         entropy=0.0\n  .rdata       entropy=0.2\n  .reloc       entropy=6.61\n  .rsrc        entropy=4.01\n  [overlay]    entropy=6.295243263244629  size=71472\n",
  "ground_truth": {
    "hash": "118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8732,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "generate random numbers using the Delphi LCG",
      "create or open file",
      "create or open registry key",
      "delay execution",
      "PEB access",
      "reference absolute stream path on Windows",
      "set file attributes",
      "create process on Windows",
      "contain loop",
      "create directory",
      "accept command line arguments",
      "read file on Windows",
      "write file on Windows",
      "query or enumerate registry value",
      "link function at runtime on Windows",
      "get geographical location",
      "set thread local storage value",
      "get thread local storage value",
      "delete file",
      "get disk size",
      "get OS version",
      "set registry value",
      "get common file path",
      "copy file",
      "persist via Run registry key",
      "contain an embedded PE file",
      "contain a thread local storage (.tls) section",
      "contain a resource (.rsrc) section",
      "compiled with Borland Delphi"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft I Service",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft I Service"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "C:\\Windows\\winhash_up.exe",
        "winhash_up.exe"
      ],
      "mutexes": [
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea.exe C:\\118a9e71a8138ddb6d196ed0a407c086e120ccdab76bc9574c648c47de5fa8ea.exe",
        "C:\\WINDOWS\\bugMAKER.bat",
        "C:\\WINDOWS\\system32\\cmd.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}