# dataset/preprocessed/12666ac4a550cb946c8df363d316f5be6b1c2c13100401315636cbf8d9635adc.json
{
  "hash": "12666ac4a550cb946c8df363d316f5be6b1c2c13100401315636cbf8d9635adc",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\12666ac4a550cb946c8df363d316f5be6b1c2c13100401315636cbf8d9635adc(null) [(null)]\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (17):\n  american animal lesbian licking  (Sarah).zip.exe  [sha256: af991d9ca46a5554...]\n  american fetish xxx uncut feet .mpeg.exe  [sha256: 4b9f1590bbd16afd...]\n  american gang bang hardcore masturbation  (Samantha).mpg.exe  [sha256: 9bd1eb599b7580fc...]\n  american horse horse girls .mpg.exe  [sha256: dcc87116cda2aaba...]\n  american nude bukkake full movie .rar.exe  [sha256: 517a9182f9ec33d7...]\n  american nude gay uncut  (Samantha).mpg.exe  [sha256: 16ac18b763476f60...]\n  beast several models feet gorgeoushorny  (Karin).avi.exe  [sha256: fe801d666596befb...]\n  black beastiality sperm [free] glans  (Ashley,Sarah).mpg.exe  [sha256: 867e97e08bab758a...]\n  blowjob full movie glans shoes .zip.exe  [sha256: bbde54e3f2c7ef5f...]\n  blowjob several models feet beautyfull .avi.exe  [sha256: 533e5c2e294d17e9...]\n  blowjob voyeur glans blondie  (Karin).mpeg.exe  [sha256: 9bc8c1b7d57469b0...]\n  brasilian handjob lingerie big  (Janette).zip.exe  [sha256: 08caed244ec88e1b...]\n  brasilian kicking bukkake [free] granny  (Anniston,Samantha).mpeg.exe  [sha256: e62f452c5f238eca...]\n  danish cumshot hardcore [free] feet .mpg.exe  [sha256: efbc30bf40d80623...]\n  danish kicking bukkake lesbian  (Janette).zip.exe  [sha256: 317a15bbe4c79582...]\n  defrag64.exe  [sha256: 7239978d5bb365c0...]\n  deinstall.exe  [sha256: 281d752565439e4d...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\blowjob several models feet beautyfull .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\lesbian [bangbus] blondie .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\japanese handjob horse big  (Curtney).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\black beastiality sperm [free] glans  (Ashley,Sarah).mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\american horse horse girls .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\danish cumshot hardcore [free] feet .mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\swedish gang bang lesbian [free] pregnant  (Jenna,Melissa).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\american gang bang hardcore masturbation  (Samantha).mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\brasilian kicking bukkake [free] granny  (Anniston,Samantha).mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\lesbian sleeping feet .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\indian action xxx masturbation  (Jade).rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\american nude gay uncut  (Samantha).mpg.exe\n  C:\\WINDOWS\\Temp\\swedish fetish sperm [bangbus] shoes .mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\american animal lesbian licking  (Sarah).zip.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\lingerie [bangbus] latex .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\brasilian handjob lingerie big  (Janette).zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\lesbian public  (Curtney).rar.exe\n  C:\\WINDOWS\\ime\\shared\\blowjob voyeur glans blondie  (Karin).mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\tyrkish cum sperm sleeping hole latex .mpg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.87998104095459  size=1285821\n",
  "ground_truth": {
    "hash": "12666ac4a550cb946c8df363d316f5be6b1c2c13100401315636cbf8d9635adc",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "american animal lesbian licking  (Sarah).zip.exe",
        "american fetish xxx uncut feet .mpeg.exe",
        "american gang bang hardcore masturbation  (Samantha).mpg.exe",
        "american horse horse girls .mpg.exe",
        "american nude bukkake full movie .rar.exe",
        "american nude gay uncut  (Samantha).mpg.exe",
        "beast several models feet gorgeoushorny  (Karin).avi.exe",
        "black beastiality sperm [free] glans  (Ashley,Sarah).mpg.exe",
        "blowjob full movie glans shoes .zip.exe",
        "blowjob several models feet beautyfull .avi.exe",
        "blowjob voyeur glans blondie  (Karin).mpeg.exe",
        "brasilian handjob lingerie big  (Janette).zip.exe",
        "brasilian kicking bukkake [free] granny  (Anniston,Samantha).mpeg.exe",
        "danish cumshot hardcore [free] feet .mpg.exe",
        "danish kicking bukkake lesbian  (Janette).zip.exe",
        "defrag64.exe",
        "deinstall.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}