# dataset/preprocessed/1422dc61e6357078f847090dffad8bdd3b760c0d8f2da26e3b8a2a53a27f5c8e.json
{
  "hash": "1422dc61e6357078f847090dffad8bdd3b760c0d8f2da26e3b8a2a53a27f5c8e",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (3):\n  ((null)) C:\\1422dc61e6357078f847090dffad8bdd3b760c0d8f2da26e3b8a2a53a27f5c8e(null) [(null)]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\kicking xxx public wifey .mpg.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african xxx [free] fishy .avi.exe  [sha256: 7749a09e740ed5ed...]\n  american xxx lingerie [free] ash sm .zip.exe  [sha256: db9d0584d4ecc010...]\n  asian nude cum girls boobs black hairunshaved .avi.exe  [sha256: a9b0350478a53751...]\n  brasilian horse girls nipples .rar.exe  [sha256: 0fbd8389b0ee55a7...]\n  british blowjob handjob [bangbus] balls .rar.exe  [sha256: d05d8069f9da18e6...]\n  cum cum [free] pregnant .mpg.exe  [sha256: 5dff831712195b88...]\n  cumshot licking feet .zip.exe  [sha256: d0ceef906ae679fb...]\n  danish gang bang girls pregnant .zip.exe  [sha256: 9d90fe8744c7507c...]\n  driver.exe  [sha256: bb9441733539f0cf...]\n  drwatson32.exe  [sha256: 8efd76022701faca...]\n  edit.exe  [sha256: dcee5c1f53499a91...]\n  explorer.exe  [sha256: acad84aac6410d53...]\n  french cum blowjob voyeur cock shoes .avi.exe  [sha256: 6cc7a60315c344e8...]\n  french gay girls glans latex .zip.exe  [sha256: 190dc594875cb3e3...]\n  game3d.exe  [sha256: 0d12dfcf96aefa8d...]\n  gay cum catfight glans  (Jenna,Sylvia).zip.exe  [sha256: 73531b5efd23decc...]\n  gay cum voyeur .zip.exe  [sha256: 1199c2b3edef5707...]\n  gay sperm public ash  (Samantha,Curtney).zip.exe  [sha256: 62df9f0a369430b1...]\n  german lingerie sperm licking .avi.exe  [sha256: de5181aacc9411a4...]\n  hardcore voyeur nipples .avi.exe  [sha256: b61bf19a2fbece68...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\danish gang bang girls pregnant .zip.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\asian nude cum girls boobs black hairunshaved .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\xxx lingerie [free] girly .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\swedish sperm sleeping .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\african xxx [free] fishy .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\hardcore voyeur nipples .avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\gay cum catfight glans  (Jenna,Sylvia).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\swedish hardcore public .mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\cum cum [free] pregnant .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\cumshot licking feet .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\french cum blowjob voyeur cock shoes .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\horse trambling several models vagina .rar.exe\n  C:\\WINDOWS\\Temp\\french gay girls glans latex .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\british blowjob handjob [bangbus] balls .rar.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\indian fucking trambling several models .mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\italian action [free] .rar.exe\n  C:\\WINDOWS\\assembly\\tmp\\porn porn lesbian sm  (Kathrin).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\german lingerie sperm licking .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\brasilian horse girls nipples .rar.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880277156829834  size=801574\n",
  "ground_truth": {
    "hash": "1422dc61e6357078f847090dffad8bdd3b760c0d8f2da26e3b8a2a53a27f5c8e",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\asian action public mature .avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\animal blowjob masturbation 50+ .zip.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\nude public .zip.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\swedish handjob girls  (Kathrin,Sarah).mpg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\horse lesbian shower  (Janette).mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\spanish beastiality licking legs  (Sonja,Anniston).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\fucking kicking voyeur cock lady .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\action beast uncut ash traffic .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\japanese handjob masturbation traffic  (Kathrin).avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\fetish [bangbus] mature  (Gina,Sarah).rar.exe",
        "C:\\Program Files (x86)\\Steam\\clientui\\images\\sharedui\\russian gay [bangbus] .zip.exe",
        "C:\\Program Files (x86)\\Steam\\controller_base\\templates\\italian lesbian sleeping .zip.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\japanese cumshot catfight nipples ash .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\danish kicking sleeping boobs .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\kicking kicking uncut glans 40+ .mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\beastiality gang bang public .avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\beast [bangbus] black hairunshaved  (Britney).rar.exe",
        "C:\\Program Files\\dotnet\\shared\\action masturbation .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "african xxx [free] fishy .avi.exe",
        "american xxx lingerie [free] ash sm .zip.exe",
        "asian nude cum girls boobs black hairunshaved .avi.exe",
        "brasilian horse girls nipples .rar.exe",
        "british blowjob handjob [bangbus] balls .rar.exe",
        "cum cum [free] pregnant .mpg.exe",
        "cumshot licking feet .zip.exe",
        "danish gang bang girls pregnant .zip.exe",
        "driver.exe",
        "drwatson32.exe",
        "edit.exe",
        "explorer.exe",
        "french cum blowjob voyeur cock shoes .avi.exe",
        "french gay girls glans latex .zip.exe",
        "game3d.exe",
        "gay cum catfight glans  (Jenna,Sylvia).zip.exe",
        "gay cum voyeur .zip.exe",
        "gay sperm public ash  (Samantha,Curtney).zip.exe",
        "german lingerie sperm licking .avi.exe",
        "hardcore voyeur nipples .avi.exe",
        "horse trambling several models vagina .rar.exe",
        "indian beastiality porn hidden cock upskirt .avi.exe",
        "indian fucking trambling several models .mpeg.exe",
        "indian horse licking gorgeoushorny .mpeg.exe",
        "install32.exe",
        "italian action [free] .rar.exe",
        "mssrv.exe",
        "nude masturbation sm  (Tatjana).mpeg.exe",
        "porn porn lesbian sm  (Kathrin).mpg.exe",
        "porn sleeping .mpg.exe",
        "screensaver.scr",
        "speedup.exe",
        "swedish cumshot hidden sweet .mpg.exe",
        "swedish hardcore public .mpg.exe",
        "swedish sperm sleeping .zip.exe",
        "undelete.exe",
        "uninstall.exe",
        "xxx lingerie [free] girly .zip.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}