# dataset/preprocessed/14f03f78cd9cea407077066adf19b308bf0efe7068c27bb7326d93d99aa1a4c9.json
{
  "hash": "14f03f78cd9cea407077066adf19b308bf0efe7068c27bb7326d93d99aa1a4c9",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\14f03f78cd9cea407077066adf19b308bf0efe7068c27bb7326d93d99aa1a4c9(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\14f03f78cd9cea40_mssrv.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american nude hardcore catfight titts mature .mpg.exe  [sha256: f254750275e9e7c3...]\n  black beastiality blowjob girls .zip.exe  [sha256: ebe0e0ea41b636b5...]\n  blowjob hot (!) .mpeg.exe  [sha256: f524929570491eaa...]\n  brasilian animal lesbian several models cock femdom .mpeg.exe  [sha256: 2741c053c5127de7...]\n  brasilian porn lingerie sleeping girly .zip.exe  [sha256: 52b4c96e84abe7be...]\n  danish fetish trambling lesbian .mpg.exe  [sha256: 92caa15a536366ab...]\n  danish gang bang trambling [milf] feet .mpeg.exe  [sha256: c867a262bef75ea9...]\n  danish horse lingerie public sweet  (Kathrin,Samantha).mpeg.exe  [sha256: 81969239eb27baab...]\n  danish nude beast [free] hole .avi.exe  [sha256: c0b92b7d46b6d3e1...]\n  danish porn bukkake voyeur hotel  (Sonja,Samantha).mpeg.exe  [sha256: e8f18d43d1fde60b...]\n  defrag.exe  [sha256: a270af92de2edcb0...]\n  defrag64.exe  [sha256: fe1825c77e87d21b...]\n  deinstall.exe  [sha256: 013d3bea59deee1c...]\n  drwatson32.exe  [sha256: 299489d87b5f6051...]\n  gay hot (!) lady .avi.exe  [sha256: db3e39a27514c2b9...]\n  german sperm [free] glans hairy .mpg.exe  [sha256: 2a14f1504f6f12a6...]\n  hardcore [bangbus] leather .mpeg.exe  [sha256: b12da1ac7635f6c7...]\n  horse sleeping hole bedroom  (Curtney).mpg.exe  [sha256: 7062a1cad8ffb79c...]\n  IE9setup.exe  [sha256: 1dcb881bb530d86f...]\n  indian fetish gay lesbian cock .mpeg.exe  [sha256: 420c437c62e3963f...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\american nude hardcore catfight titts mature .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\indian gang bang xxx [bangbus] latex .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\indian fetish gay lesbian cock .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\brasilian animal lesbian several models cock femdom .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\xxx [milf] feet black hairunshaved .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\horse sleeping hole bedroom  (Curtney).mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\danish nude beast [free] hole .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\german sperm [free] glans hairy .mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\tyrkish nude lingerie uncut .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\tyrkish animal horse catfight cock boots .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\xxx [bangbus] .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\danish porn bukkake voyeur hotel  (Sonja,Samantha).mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\danish fetish trambling lesbian .mpg.exe\n  C:\\WINDOWS\\Temp\\brasilian porn lingerie sleeping girly .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\blowjob hot (!) .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\indian handjob lesbian lesbian glans ejaculation .mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\trambling lesbian  (Liz).zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\tyrkish cum lingerie masturbation sweet .avi.exe\n  C:\\WINDOWS\\ime\\shared\\tyrkish action beast [free] .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.8805694580078125  size=1213301\n",
  "ground_truth": {
    "hash": "14f03f78cd9cea407077066adf19b308bf0efe7068c27bb7326d93d99aa1a4c9",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\kicking fetish uncut .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\blowjob beastiality hot (!) .mpeg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\african horse beastiality several models feet hotel  (Janette).mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\malaysia lesbian lesbian [bangbus] beautyfull .mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\malaysia action beast several models cock granny .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\african xxx gang bang public 50+ .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american nude hardcore catfight titts mature .mpg.exe",
        "black beastiality blowjob girls .zip.exe",
        "blowjob hot (!) .mpeg.exe",
        "brasilian animal lesbian several models cock femdom .mpeg.exe",
        "brasilian porn lingerie sleeping girly .zip.exe",
        "danish fetish trambling lesbian .mpg.exe",
        "danish gang bang trambling [milf] feet .mpeg.exe",
        "danish horse lingerie public sweet  (Kathrin,Samantha).mpeg.exe",
        "danish nude beast [free] hole .avi.exe",
        "danish porn bukkake voyeur hotel  (Sonja,Samantha).mpeg.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "drwatson32.exe",
        "gay hot (!) lady .avi.exe",
        "german sperm [free] glans hairy .mpg.exe",
        "hardcore [bangbus] leather .mpeg.exe",
        "horse sleeping hole bedroom  (Curtney).mpg.exe",
        "indian fetish gay lesbian cock .mpeg.exe",
        "indian gang bang xxx [bangbus] latex .avi.exe",
        "indian handjob lesbian lesbian glans ejaculation .mpg.exe",
        "italian porn lingerie licking sweet .avi.exe",
        "japanese porn beast big cock latex  (Liz).rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "russian nude xxx [bangbus]  (Tatjana).zip.exe",
        "screensaver.scr",
        "setup.exe",
        "soccer.exe",
        "start.exe",
        "trambling lesbian  (Liz).zip.exe",
        "tyrkish action beast [free] .avi.exe",
        "tyrkish animal horse catfight cock boots .rar.exe",
        "tyrkish cum lingerie masturbation sweet .avi.exe",
        "tyrkish nude lingerie uncut .mpg.exe",
        "xxx [bangbus] .avi.exe",
        "xxx [milf] feet black hairunshaved .mpeg.exe",
        "xxx uncut femdom .zip.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}