# dataset/preprocessed/159cae7e21c9456328a7190ee3967c9cf44f27fcba9c9ec1e7344e0ea963c0b2.json
{
  "hash": "159cae7e21c9456328a7190ee3967c9cf44f27fcba9c9ec1e7344e0ea963c0b2",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\159cae7e21c9456328a7190ee3967c9cf44f27fcba9c9ec1e7344e0ea963c0b2(null) [(null)]\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american horse bukkake full movie gorgeoushorny .avi.exe  [sha256: b94b5e87d7e6570b...]\n  black cum bukkake sleeping high heels .mpg.exe  [sha256: 3f405a6160c8f0f3...]\n  black gang bang lesbian hot (!) .mpg.exe  [sha256: fafc10f48bc6f039...]\n  black porn lesbian [bangbus]  (Janette).rar.exe  [sha256: 3101d3041169fa8a...]\n  blowjob [bangbus] .avi.exe  [sha256: 58b8e82c3336804f...]\n  blowjob [milf] .zip.exe  [sha256: e6b6c70099948e5b...]\n  brasilian gang bang sperm lesbian black hairunshaved  (Christine,Karin).zip.exe  [sha256: 915399ef8db46b20...]\n  danish action sperm public 40+ .zip.exe  [sha256: 4bec97706ec446a8...]\n  danish gang bang lingerie catfight cock .mpg.exe  [sha256: a950d26b05f6c197...]\n  danish gang bang sperm [free]  (Liz).mpeg.exe  [sha256: 0be69f4be35dc936...]\n  defrag64.exe  [sha256: 0cb46ae324a0cf0f...]\n  drwatson32.exe  [sha256: b7adff814fdef37e...]\n  freeftp.exe  [sha256: ca0a4c0af6fa0f33...]\n  fucking [milf] black hairunshaved .mpg.exe  [sha256: 8ee753b22208b340...]\n  gang bang gay masturbation .mpg.exe  [sha256: 9c987c10f57280b0...]\n  gay [bangbus] castration .avi.exe  [sha256: 9032ff18602ee258...]\n  hardcore masturbation  (Janette).mpeg.exe  [sha256: 9c62a54992881e15...]\n  hardcore voyeur feet .zip.exe  [sha256: 351527f16bcbd9d2...]\n  horse hidden  (Tatjana).mpeg.exe  [sha256: 8439b5d579518fa5...]\n  horse voyeur blondie .mpeg.exe  [sha256: d1575c1f1781fb0a...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\trambling sleeping glans .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\black cum bukkake sleeping high heels .mpg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\horse voyeur blondie .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\gay [bangbus] castration .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\horse hidden  (Tatjana).mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\black gang bang lesbian hot (!) .mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\blowjob [milf] .zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\american horse bukkake full movie gorgeoushorny .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\brasilian gang bang sperm lesbian black hairunshaved  (Christine,Karin).zip.exe\n  C:\\Program Files\\Movie Maker\\Shared\\xxx uncut .avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\swedish animal bukkake licking hole high heels .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\danish action sperm public 40+ .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\danish gang bang lingerie catfight cock .mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\gang bang gay masturbation .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\italian action horse catfight glans .mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\japanese beastiality bukkake lesbian  (Karin).rar.exe\n  C:\\WINDOWS\\assembly\\tmp\\fucking [milf] black hairunshaved .mpg.exe\n  C:\\WINDOWS\\ime\\shared\\xxx full movie fishy .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\hardcore masturbation  (Janette).mpeg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.879977226257324  size=402682\n",
  "ground_truth": {
    "hash": "159cae7e21c9456328a7190ee3967c9cf44f27fcba9c9ec1e7344e0ea963c0b2",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "IE9setup.exe",
        "american horse bukkake full movie gorgeoushorny .avi.exe",
        "black cum bukkake sleeping high heels .mpg.exe",
        "black gang bang lesbian hot (!) .mpg.exe",
        "black porn lesbian [bangbus]  (Janette).rar.exe",
        "blowjob [bangbus] .avi.exe",
        "blowjob [milf] .zip.exe",
        "brasilian gang bang sperm lesbian black hairunshaved  (Christine,Karin).zip.exe",
        "danish action sperm public 40+ .zip.exe",
        "danish gang bang lingerie catfight cock .mpg.exe",
        "danish gang bang sperm [free]  (Liz).mpeg.exe",
        "defrag64.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking [milf] black hairunshaved .mpg.exe",
        "gang bang gay masturbation .mpg.exe",
        "gay [bangbus] castration .avi.exe",
        "hardcore masturbation  (Janette).mpeg.exe",
        "hardcore voyeur feet .zip.exe",
        "horse hidden  (Tatjana).mpeg.exe",
        "horse voyeur blondie .mpeg.exe",
        "italian action horse catfight glans .mpeg.exe",
        "italian kicking xxx public black hairunshaved .mpeg.exe",
        "japanese action lesbian masturbation feet granny .avi.exe",
        "japanese beastiality bukkake lesbian  (Karin).rar.exe",
        "lesbian masturbation hotel .rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "soccer.exe",
        "start.exe",
        "swedish animal bukkake licking hole high heels .zip.exe",
        "trambling sleeping glans .mpeg.exe",
        "tyrkish kicking hardcore licking .rar.exe",
        "undelete.exe",
        "xxx full movie fishy .zip.exe",
        "xxx masturbation feet granny  (Sarah).rar.exe",
        "xxx uncut .avi.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}