# dataset/preprocessed/15a889bdeda88c96a1dbdf95ff5ebd0be90a235edf10a814c51ad623802316e3.json
{
  "hash": "15a889bdeda88c96a1dbdf95ff5ebd0be90a235edf10a814c51ad623802316e3",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  3d5d9f6c2c1dd7e8ff38a631da84c246.exe\n    C:\\Users\\<USER>\\Downloads\\3d5d9f6c2c1dd7e8ff38a631da84c246.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\15a889bdeda88c96a1dbdf95ff5ebd0be90a235edf10a814c51ad623802316e3(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\3d5d9f6c2c1dd7e8ff38a631da84c246.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (19):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n  - packed with UPX\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african cumshot fetish masturbation titts  (Melissa,Sandy).rar.exe  [sha256: ee81c77ba9c249f6...]\n  beast fucking public mature .mpeg.exe  [sha256: e05dbaa37b90dd12...]\n  beastiality big latex .avi.exe  [sha256: bfc5d63a97dc14cd...]\n  blowjob hidden vagina sweet  (Britney,Kathrin).avi.exe  [sha256: 0d18845a0199884e...]\n  brasilian cum uncut ash beautyfull  (Melissa,Tatjana).mpg.exe  [sha256: e6b6ba3dd9390a6a...]\n  british cumshot voyeur black hairunshaved  (Kathrin).mpg.exe  [sha256: fbbbbd3765d9721c...]\n  bukkake public boots  (Karin).rar.exe  [sha256: 032949d7751d4922...]\n  canadian action action uncut 50+ .avi.exe  [sha256: 9aa992c3a8d5da05...]\n  danish action fetish [free] hairy  (Sandy,Janette).rar.exe  [sha256: 0595ea2e3df5d545...]\n  danish cum lesbian .avi.exe  [sha256: 77d3023cc079d3f8...]\n  danish handjob several models girly .mpg.exe  [sha256: 918c7f21a4a0adef...]\n  defrag.exe  [sha256: b7b04e3701bc8972...]\n  defrag64.exe  [sha256: 5278344be4aaacc5...]\n  deinstall.exe  [sha256: 45e8b4237185f5cb...]\n  driver.exe  [sha256: 7f76fdd200b70f71...]\n  fetish several models .mpg.exe  [sha256: 802dc2f88a0f104f...]\n  freeftp.exe  [sha256: 33dc52953e351cee...]\n  game3d.exe  [sha256: bfd35cb90c260bc4...]\n  german action bukkake public nipples  (Sonja).avi.exe  [sha256: 7aba5131c5550095...]\n  german beast licking titts .rar.exe  [sha256: 8cba87afb4f1fd71...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\british cumshot voyeur black hairunshaved  (Kathrin).mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\malaysia trambling beastiality public traffic  (Liz,Sonja).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\italian bukkake public .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\danish cum lesbian .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\lesbian voyeur blondie .rar.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\beast fucking public mature .mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\blowjob hidden vagina sweet  (Britney,Kathrin).avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\fetish several models .mpg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\german action bukkake public nipples  (Sonja).avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\russian beast gay [milf] boobs .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\lingerie [free] leather .mpeg.exe\n  C:\\WINDOWS\\Temp\\indian beast fetish catfight .avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\brasilian cum uncut ash beautyfull  (Melissa,Tatjana).mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\bukkake public boots  (Karin).rar.exe\n  C:\\WINDOWS\\assembly\\temp\\canadian action action uncut 50+ .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\danish handjob several models girly .mpg.exe\n  C:\\WINDOWS\\ime\\shared\\japanese hardcore fetish big latex .mpg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\xxx licking black hairunshaved .avi.exe\n  C:\\WINDOWS\\security\\templates\\porn big gorgeoushorny .avi.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  UPX0         entropy=6.29\n  UPX1         entropy=4.91\n  UPX2         entropy=3.31\n  .imports     entropy=4.08\n  [overlay]    entropy=7.88012170791626  size=1634420\n",
  "ground_truth": {
    "hash": "15a889bdeda88c96a1dbdf95ff5ebd0be90a235edf10a814c51ad623802316e3",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8732,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "packed with UPX"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "packed with UPX",
      "(internal) packer file limitation",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\fucking several models swallow .zip.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\kicking blowjob catfight mature  (Tatjana).rar.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\gang bang cumshot licking high heels  (Karin).mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\french cum beastiality catfight .zip.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\japanese handjob [free] sm  (Ashley,Sonja).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\african cumshot sleeping nipples 40+ .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\indian xxx animal uncut hairy .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\horse voyeur .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\russian nude nude several models  (Melissa,Samantha).mpg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\gay masturbation balls .zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\american blowjob porn big ash .avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\danish fucking cumshot hot (!) nipples .rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\porn trambling several models titts latex  (Melissa,Sandy).avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\nude gay uncut beautyfull  (Kathrin,Jenna).zip.exe",
        "C:\\Program Files\\dotnet\\shared\\black porn xxx catfight glans .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "african cumshot fetish masturbation titts  (Melissa,Sandy).rar.exe",
        "beast fucking public mature .mpeg.exe",
        "beastiality big latex .avi.exe",
        "blowjob hidden vagina sweet  (Britney,Kathrin).avi.exe",
        "brasilian cum uncut ash beautyfull  (Melissa,Tatjana).mpg.exe",
        "british cumshot voyeur black hairunshaved  (Kathrin).mpg.exe",
        "bukkake public boots  (Karin).rar.exe",
        "canadian action action uncut 50+ .avi.exe",
        "danish action fetish [free] hairy  (Sandy,Janette).rar.exe",
        "danish cum lesbian .avi.exe",
        "danish handjob several models girly .mpg.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "fetish several models .mpg.exe",
        "freeftp.exe",
        "game3d.exe",
        "german action bukkake public nipples  (Sonja).avi.exe",
        "german beast licking titts .rar.exe",
        "indian beast fetish catfight .avi.exe",
        "italian bukkake public .mpg.exe",
        "italian handjob cumshot voyeur pregnant .mpeg.exe",
        "japanese hardcore fetish big latex .mpg.exe",
        "lesbian voyeur blondie .rar.exe",
        "lingerie [free] leather .mpeg.exe",
        "malaysia gang bang horse [bangbus] sm  (Christine).mpg.exe",
        "malaysia trambling beastiality public traffic  (Liz,Sonja).mpeg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}