# dataset/preprocessed/1617f209390c6b4294ff2280817b9e6c9cfa2681e8aa801a6e0539bda0efe4a7.json
{
  "hash": "1617f209390c6b4294ff2280817b9e6c9cfa2681e8aa801a6e0539bda0efe4a7",
  "report_text": "--- NETWORK ACTIVITY ---\n\nSERVICES OPENED (1):\n  Messenger\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (2):\n  C:\\Users\\<USER>\\Downloads\\90707182ccd97732a626a8ce57d0681f.virus.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (19):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n  - packed with UPX\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  C:\\Windows\\mssrv.exe  [sha256: 9620af72c54469ff...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\american kicking fucking licking  (Janette).mpeg.exe  [sha256: 7f4b0474e1fd9258...]\n  C:\\Program Files\\dotnet\\shared\\japanese kicking lingerie hot (!) mistress  (Christine,Sarah).rar.exe  [sha256: aa2370d846ad564a...]\n  C:\\Program Files\\DVD Maker\\Shared\\swedish cum lingerie catfight glans upskirt  (Curtney).zip.exe  [sha256: a97b7394d710d274...]\n  C:\\Program Files\\Windows Journal\\Templates\\russian beastiality blowjob sleeping feet bedroom  (Janette).avi.exe  [sha256: 97bbd80e33b8d471...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\black cumshot bukkake licking latex .avi.exe  [sha256: 02c0e42574df5799...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\japanese gang bang bukkake voyeur balls .mpg.exe  [sha256: 1adf96c41457ca95...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\bukkake big hole .zip.exe  [sha256: e83a925f6904573a...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\italian cumshot gay sleeping .rar.exe  [sha256: f5d4b18854496f77...]\n  C:\\Program Files (x86)\\dotnet\\shared\\danish horse sperm [bangbus] glans .avi.exe  [sha256: 9f5b97937c66fc15...]\n  C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\danish cum lingerie several models granny .mpeg.exe  [sha256: 0e0faed3a1859710...]\n  C:\\Program Files (x86)\\Google\\Update\\Download\\xxx licking upskirt .mpeg.exe  [sha256: 80c3fb7e97c3ed19...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\fucking voyeur hairy  (Sonja,Melissa).mpg.exe  [sha256: 4101a9a114d81387...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\swedish nude xxx sleeping  (Karin).mpeg.exe  [sha256: 6e685a24124fd536...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\swedish gang bang blowjob uncut hole  (Kathrin,Jade).avi.exe  [sha256: b5062cecbb48d661...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\japanese handjob gay big feet ejaculation  (Curtney).avi.exe  [sha256: 28af3b63ae7f0e28...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\japanese nude fucking girls redhair .mpg.exe  [sha256: 8f5888b494025a6c...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\italian beastiality xxx [bangbus] traffic  (Ashley,Sylvia).rar.exe  [sha256: 0c4677cb6509ec05...]\n  C:\\Program Files\\dotnet\\shared\\russian fetish lingerie hidden swallow .avi.exe  [sha256: 1267a3d9742e1bab...]\n  C:\\Program Files\\DVD Maker\\Shared\\lesbian licking .avi.exe  [sha256: 57f629a8711ad2a8...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\Administrator\\Application Data\\Macromedia\\Flash Player\\#SharedObjects\\indian blowjob several models .rar.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Bin\\bse_temp\\asian cum [free] .avi.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Extension\\com.sogou.share\\tyrkish blowjob horse several models 50+ .zip.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Temp\\gay beastiality [free] granny .zip.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQPCMgr\\Download\\horse horse full movie balls .rar.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\Download\\swedish fucking hot (!) castration .avi.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\action cum lesbian cock  (Ashley,Christine).rar.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\spanish bukkake xxx big  (Sylvia,Gina).rar.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\Temp\\indian gay public .zip.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\TXSSO\\SSOTemp\\swedish porn gang bang [bangbus]  (Liz,Kathrin).mpg.exe\n  C:\\Documents and Settings\\Administrator\\CMB\\PB40\\Temp\\asian gay several models high heels .avi.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\french fetish big girly .avi.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temporary Internet Files\\swedish trambling [milf] feet .rar.exe\n  C:\\Documents and Settings\\Administrator\\Templates\\german horse full movie  (Anniston).avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\american handjob public gorgeoushorny  (Christine,Curtney).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\hardcore lesbian .rar.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\nude masturbation castration .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temporary Internet Files\\italian cum big .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\danish trambling public .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\swedish porn lingerie voyeur .mpeg.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  UPX0         entropy=6.29\n  UPX1         entropy=4.91\n  UPX2         entropy=3.31\n  .imports     entropy=4.08\n  [overlay]    entropy=7.879419803619385  size=331261\n",
  "ground_truth": {
    "hash": "1617f209390c6b4294ff2280817b9e6c9cfa2681e8aa801a6e0539bda0efe4a7",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8732,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "packed with UPX"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "packed with UPX",
      "(internal) packer file limitation",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\japanese gang bang bukkake voyeur balls .mpg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\tyrkish kicking beast voyeur glans redhair .avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\bukkake big hole .zip.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\fucking hidden blondie  (Anniston,Samantha).mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\italian cumshot gay sleeping .rar.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\tyrkish gang bang lesbian licking hole .mpg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\american horse fucking lesbian bondage .mpg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\danish cum lingerie several models granny .mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\danish cum blowjob several models cock .mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\xxx licking upskirt .mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\fucking voyeur hairy  (Sonja,Melissa).mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\gay hot (!) .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\swedish nude xxx sleeping  (Karin).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\tyrkish animal gay [milf] cock high heels  (Karin).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\japanese handjob gay big feet ejaculation  (Curtney).avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\swedish cum bukkake licking leather .mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\russian action hardcore hot (!) feet ejaculation  (Samantha).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\swedish gang bang blowjob uncut hole  (Kathrin,Jade).avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\bukkake [free] redhair .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\japanese nude fucking girls redhair .mpg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\danish horse sperm [bangbus] glans .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\russian kicking trambling several models cock .zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\american kicking fucking licking  (Janette).mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\italian beastiality xxx [bangbus] traffic  (Ashley,Sylvia).rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\lesbian licking .avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\swedish cum lingerie catfight glans upskirt  (Curtney).zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\japanese action horse licking titts bedroom .rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\russian beastiality blowjob sleeping feet bedroom  (Janette).avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\black cumshot bukkake licking latex .avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\bukkake public cock sweet .mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\japanese kicking lingerie hot (!) mistress  (Christine,Sarah).rar.exe",
        "C:\\Program Files\\dotnet\\shared\\russian fetish lingerie hidden swallow .avi.exe",
        "C:\\Windows\\mssrv.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Jujubox"
  ]
}