# dataset/preprocessed/16403e86d1d26fdf613358b9f382d483a0bb8e0359f5ac03b3a047f78a9ad8fb.json
{
  "hash": "16403e86d1d26fdf613358b9f382d483a0bb8e0359f5ac03b3a047f78a9ad8fb",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  danish-sperm-public.pe32\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n\nCOMMAND EXECUTIONS (3):\n  ((null)) C:\\16403e86d1d26fdf613358b9f382d483a0bb8e0359f5ac03b3a047f78a9ad8fb(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\16403e86d1d26fdf_danish sperm public .mpeg.exe\n  \"C:\\Windows\\system32\\rundll32.exe\" C:\\Windows\\system32\\shell32.dll,OpenAs_RunDLL C:\\Users\\Petra\\AppData\\Local\\Temp\\danish-sperm-public.pe32\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american horse beast hot (!)  (Karin).rar.exe  [sha256: e809e554fb3e2d5d...]\n  beast lesbian .avi.exe  [sha256: 66c69882a5aa6f2e...]\n  black cumshot horse big cock traffic .mpeg.exe  [sha256: 86ad741c825161a1...]\n  black fetish gay several models  (Liz).mpeg.exe  [sha256: 8799698611c31d1c...]\n  black horse xxx voyeur hole young .avi.exe  [sha256: afa8b9097b06d48f...]\n  black nude xxx licking cock fishy  (Liz).rar.exe  [sha256: 2cb3d098d98eb7f7...]\n  blowjob catfight glans latex  (Sarah).zip.exe  [sha256: e821aebea735a659...]\n  brasilian gang bang hardcore [milf] blondie .avi.exe  [sha256: fa428a5012ec5f14...]\n  bukkake hidden .mpeg.exe  [sha256: 629f498bf3de7d3c...]\n  bukkake [free] shoes .avi.exe  [sha256: 93747d2fe7290eed...]\n  danish handjob sperm lesbian .zip.exe  [sha256: 2fb4f78c6203eef7...]\n  danish kicking gay big feet wifey .zip.exe  [sha256: 93f294d2e408b966...]\n  defrag.exe  [sha256: 09455b4360b4d7c2...]\n  defrag64.exe  [sha256: 32b065dad84e7b17...]\n  gay big pregnant .mpeg.exe  [sha256: 9c8c7a598ce31642...]\n  gay hidden beautyfull .rar.exe  [sha256: 04aafa8dfde20123...]\n  gay voyeur sweet .avi.exe  [sha256: 6c6698ab0a19fb21...]\n  gay [free] glans .zip.exe  [sha256: 6b531d9d9703684a...]\n  hardcore catfight wifey .avi.exe  [sha256: 59ec0dbf8b6b7185...]\n  horse [bangbus]  (Tatjana).avi.exe  [sha256: 2c658d77336dd98e...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\american horse beast hot (!)  (Karin).rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\italian cum horse hidden .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\bukkake hidden .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\bukkake [free] shoes .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\gay [free] glans .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\black horse xxx voyeur hole young .avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\black cumshot horse big cock traffic .mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\danish handjob sperm lesbian .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\black nude xxx licking cock fishy  (Liz).rar.exe\n  C:\\Program Files\\Movie Maker\\Shared\\gay hidden beautyfull .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\gay voyeur sweet .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\brasilian gang bang hardcore [milf] blondie .avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\trambling [milf] titts gorgeoushorny .mpg.exe\n  C:\\WINDOWS\\Temp\\trambling several models lady .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\japanese action lesbian masturbation glans .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\beast lesbian .avi.exe\n  C:\\WINDOWS\\assembly\\temp\\black fetish gay several models  (Liz).mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\italian action blowjob voyeur swallow .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\swedish nude hardcore voyeur titts .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\MACHINE\\HexStream = 53004F004600540057004100520045005C004D006900630072006F0073006F00660074005C00570069006E0064\n\nMUTEXES CREATED (1):\n  mutex666\n\nPROCESSES CREATED (1):\n  rundll32.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.8801164627075195  size=836504\n",
  "ground_truth": {
    "hash": "16403e86d1d26fdf613358b9f382d483a0bb8e0359f5ac03b3a047f78a9ad8fb",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\swedish animal porn [milf] gorgeoushorny .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\british sperm nude hot (!)  (Liz,Sonja).mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\spanish porn porn masturbation balls .avi.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\chinese gang bang horse licking stockings .zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\kicking cum hidden nipples blondie .mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\asian beastiality hidden ash .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "\\Program Files\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\lingerie big black hairunshaved  (Jenna,Kathrin).mpg.exe",
        "\\Program Files\\DVD Maker\\Shared\\french kicking licking girly .rar.exe",
        "\\Program Files\\Google\\Update\\1.3.34.11\\Recovery\\GUR4F00.tmp\\african horse full movie mistress .zip.exe",
        "\\Program Files\\Google\\Update\\Download\\italian action handjob girls nipples .mpg.exe",
        "\\Windows\\mssrv.exe",
        "american horse beast hot (!)  (Karin).rar.exe",
        "beast lesbian .avi.exe",
        "black cumshot horse big cock traffic .mpeg.exe",
        "black fetish gay several models  (Liz).mpeg.exe",
        "black horse xxx voyeur hole young .avi.exe",
        "black nude xxx licking cock fishy  (Liz).rar.exe",
        "blowjob catfight glans latex  (Sarah).zip.exe",
        "brasilian gang bang hardcore [milf] blondie .avi.exe",
        "bukkake [free] shoes .avi.exe",
        "bukkake hidden .mpeg.exe",
        "danish handjob sperm lesbian .zip.exe",
        "danish kicking gay big feet wifey .zip.exe",
        "defrag.exe",
        "defrag64.exe",
        "gay [free] glans .zip.exe",
        "gay big pregnant .mpeg.exe",
        "gay hidden beautyfull .rar.exe",
        "gay voyeur sweet .avi.exe",
        "hardcore catfight wifey .avi.exe",
        "horse [bangbus]  (Tatjana).avi.exe",
        "italian action blowjob voyeur swallow .mpeg.exe",
        "italian cum horse hidden .avi.exe",
        "italian porn xxx several models .avi.exe",
        "japanese action lesbian masturbation glans .mpeg.exe",
        "lesbian voyeur hole .avi.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "screensaver.scr",
        "setup.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": [
        "rundll32.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "SNDBOX",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}