# dataset/preprocessed/165ee30b935c9aca323363db01927b9dac88b7fa32e230f88c54f3daee5b158b.json
{
  "hash": "165ee30b935c9aca323363db01927b9dac88b7fa32e230f88c54f3daee5b158b",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (5):\n  elegan_786444.el.funpic.org\n  elegan_786444.el.funpic.org \u2192 \n  ftp.tripod.com \u2192 209.202.252.54\n  griptoloji.host-ed.net\n  griptoloji.host-ed.net \u2192 \n\nIP TRAFFIC (1):\n  209.202.252.54:21 (TCP)\n\nSERVICES OPENED (2):\n  RASMAN\n  Sens\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  165a4aa97cdb46b22d7d8e456edbc72d.exe\n    73b12ce12dea223e3f12e0e19b669c54.exe\n  4f3093e2aa57962314245e04ad6b6cf2.exe\n    d1c56374fff0243832b8696d133b7861.exe\n    d1c56374fff0243832b8696d133b7861.exe\n      7c34a78b4422b61608ab858fb3833bc4.exe\n      1c799b01c13133d76f5dc410be4e41a1.exe\n      1c799b01c13133d76f5dc410be4e41a1.exe\n    d1c56374fff0243832b8696d133b7861.exe\n    d1c56374fff0243832b8696d133b7861.exe\n    0f092e237d6e31a4a6939b61477d8d2d.exe\n    d1c56374fff0243832b8696d133b7861.exe\n\nCOMMAND EXECUTIONS (7):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n  C:\\Program Files (x86)\\a4685e47\\jusched.exe \n  C:\\Windows\\system32\\DllHost.exe /Processid:{F9717507-6651-4EDB-BFF7-AE615179BCCF}\n  C:\\Windows\\system32\\wbem\\wmiprvse.exe -Embedding\n  C:\\Windows\\System32\\slui.exe -Embedding\n  C:\\Windows\\system32\\svchost.exe -k netsvcs\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - parse PE header\n  - enumerate files on Windows\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - parse PE header\n  - Discovery::File and Directory Discovery [E1083]\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - debug build\n  - contain a resource (.rsrc) section\n  - contains PDB path\n\nHIGHLIGHTED CALLS:\n  SetFileTime\n  C:\\Program Files (x86)\\a4685e47\\jusched.exe\n  C:\\Windows\\system32\\DllHost.exe /Processid:{F9717507-6651-4EDB-BFF7-AE615179BCCF}\n  C:\\Windows\\system32\\wbem\\wmiprvse.exe -Embedding\n  C:\\Windows\\System32\\slui.exe -Embedding\n  C:\\Windows\\system32\\svchost.exe -k netsvcs\n  kernel32.dll.FreeLibrary\n  kernel32.dll.FindFirstFileA\n  kernel32.dll.WriteFile\n  kernel32.dll.ReadFile\n  shell32.dll.ShellExecuteA\n  wininet.dll.InternetConnectA\n  samlib.dll.SamConnect\n  kernel32.dll.RegOpenKeyExW\n  kernel32.dll.IsDebuggerPresent\n\n--- ARTIFACTS ---\n\nFILES DROPPED (2):\n  jusched.exe  [sha256: 37727f9d31254d19...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: 40c894dcf19e9d77...]\n\nFILES WRITTEN (26):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files (x86)\\a4685e47\\a4685e47\n  C:\\Program Files (x86)\\a4685e47\\jusched.exe\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Microsoft\\Windows\\Temporary Internet Files\\counters.dat\n  C:\\WINDOWS\\Tasks\\Update23.job\n  C:\\Windows\\ServiceProfiles\\NetworkService\\AppData\\Roaming\\Microsoft\\SoftwareProtectionPlatform\\Cache\\cache.dat\n  C:\\Windows\\SoftwareDistribution\\DataStore\\DataStore.edb\n  C:\\Windows\\SoftwareDistribution\\DataStore\\Logs\\edb.chk\n  C:\\Windows\\Tasks\\Update23.job\n  C:\\Windows\\Temp\\fwtsqmfile03.sqm\n  C:\\Windows\\appcompat\\Programs\\RecentFileCache.bcf\n  C:\\Windows\\sysnative\\Tasks\\Update23\n  C:\\Windows\\sysnative\\wbem\\repository\\INDEX.BTR\n  C:\\Windows\\sysnative\\wbem\\repository\\MAPPING1.MAP\n  C:\\Windows\\sysnative\\wbem\\repository\\MAPPING2.MAP\n  C:\\Windows\\sysnative\\wbem\\repository\\MAPPING3.MAP\n  C:\\Windows\\sysnative\\wbem\\repository\\OBJECTS.DATA\n  C:\\Windows\\sysnative\\wbem\\repository\\WRITABLE.TST\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:JavaUpdate23\n\nOTHER REGISTRY KEYS SET (2):\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\ProxyEnable\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\\n\nMUTEXES CREATED (2):\n  RasPbFile\n  ShimCacheMutex\n\nMUTEXES OPENED (3):\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!history!history.ie5!\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  Local\\c:!users!<USER>!appdata!roaming!microsoft!windows!cookies!\n\nPROCESSES CREATED (1):\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  [Overlay] PDB 2.0 file link\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.88\n  .rdata       entropy=4.37\n  .data        entropy=0.44\n  .idata       entropy=2.84\n  .rsrc        entropy=2.25\n  .reloc       entropy=7.83\n  [overlay]    entropy=7.70659875869751  size=17536\n",
  "ground_truth": {
    "hash": "165ee30b935c9aca323363db01927b9dac88b7fa32e230f88c54f3daee5b158b",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9265,
    "malicious_count": 63,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "parse PE header",
      "enumerate files on Windows",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "debug build",
      "contain a resource (.rsrc) section",
      "contains PDB path"
    ],
    "iocs": {
      "ips": [
        "209.202.252.54"
      ],
      "domains": [
        "elegan_786444.el.funpic.org",
        "ftp.tripod.com",
        "griptoloji.host-ed.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "jusched.exe"
      ],
      "mutexes": [
        "RasPbFile",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \""
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox",
    "Yomi Hunter"
  ]
}