# dataset/preprocessed/1660d3fd61e177f897702329112538de2bb7823bc9650d75f92db5cc96cc123c.json
{
  "hash": "1660d3fd61e177f897702329112538de2bb7823bc9650d75f92db5cc96cc123c",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  CSQOVP.EXE.exe\n    C:\\Program Files (x86)\\EDZES.EXE\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (28):\n  - set file attributes\n  - encode data using XOR\n  - access PEB ldr_data\n  - hide graphical window\n  - accept command line arguments\n  - query or enumerate registry value\n  - link function at runtime on Windows\n  - get geographical location\n  - get file size\n  - find graphical window\n  - get disk size\n  - get system information on Windows\n  - read clipboard data\n  - capture screenshot\n  - log keystrokes via polling\n  - link many functions at runtime\n  - log keystrokes\n  - get keyboard layout\n  - enumerate gui resources\n  - encode data using Base64\n  - create service\n  - delete service\n  - persist via Windows service\n  - get hostname\n  - get disk information\n  - enumerate files on Windows\n  - packed with Mpress\n  - reference Base64 string\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Process::Create Thread [C0038]\n  - File System::Create File [C0016]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - Process::Resume Thread [C0054]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - access PEB ldr_data\n  - hide graphical window\n  - set global application hook\n  - inspect load icon resource\n  - Communication::Socket Communication::Create TCP Socket [C0001.011]\n  - Process::Check Mutex [C0043]\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - contain loop\n  - Execution::Command and Scripting Interpreter [E1059]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - link function at runtime on Windows\n\n--- ARTIFACTS ---\n\nFILES DROPPED (13):\n  HJUV.EXE  [sha256: 09422bab4023a615...]\n  LFQRGF.EXE  [sha256: ef0f319f461f8bfc...]\n  TPI.EXE  [sha256: 952f96eff5103ce0...]\n  WJF.EXE  [sha256: e27b12acf98b4edb...]\n  YHDP.EXE  [sha256: 659b67c021b6d37e...]\n  ZZB.EXE  [sha256: 06f5c9143400365e...]\n  C:\\Program Files (x86)\\EDZES.EXE  [sha256: fc40a692fd24405e...]\n  C:\\Program Files (x86)\\JDHDD.EXE  [sha256: 154856df6578baa4...]\n  C:\\Program Files\\NCBE.EXE  [sha256: dba8f5d925ede2ed...]\n  C:\\Program Files\\UQMH.EXE  [sha256: 14b7be417db24ef1...]\n  C:\\PerfLogs\\NFQQ.EXE  [sha256: cf1cfc36a1f1a0b0...]\n  C:\\PerfLogs\\FXLM.EXE  [sha256: bf3b9cfc5dbbff3e...]\n  C:\\$Recycle.Bin\\ECR.EXE  [sha256: e6ec4671f477b8f5...]\n\nFILES WRITTEN (10):\n  C:\\Documents and Settings\\HJUV.EXE\n  C:\\Documents and Settings\\YHDP.EXE\n  C:\\RECYCLER\\LFQRGF.EXE\n  C:\\System Volume Information\\EPTLNT.EXE\n  C:\\System Volume Information\\HJGYPM.EXE\n  C:\\System Volume Information\\YWHXFY.EXE\n  C:\\cuckoo\\TPI.EXE\n  C:\\cuckoo\\WJF.EXE\n  C:\\cuckoo\\ZZB.EXE\n  c:\\filedebug\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\JVXUB.EXE = C:\\$Recycle.Bin\\ECR.EXE\n\nOTHER REGISTRY KEYS SET (6):\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Setup\\workfile\n  HKCR\\txtfile\\shell\\open\\command\\\n  HKCR\\QQQ.file\\shell\\open\\command\\\n  HKCR\\QQQfile\\shell\\open\\command\\\n  HKCR\\inffile\\shell\\open\\command\\\n  HKCR\\exefile\\shell\\open\\command\\\n\nMUTEXES CREATED (2):\n  ShimCacheMutex\n  qqpass7\n\nPROCESSES CREATED (1):\n  C:\\cuckoo\\ZZB.EXE\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc, WinExec\n\nCOMPILER AND PACKER:\n  [Packer] MPRESS\n  [Compiler] Borland Delphi\n\nPE SECTIONS (entropy):\n  .MPRESS1     entropy=5.88\n  .MPRESS2     entropy=5.85\n  .rsrc        entropy=2.49\n  .imports     entropy=4.54\n  [overlay]    entropy=0.31630972027778625  size=41069\n",
  "ground_truth": {
    "hash": "1660d3fd61e177f897702329112538de2bb7823bc9650d75f92db5cc96cc123c",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9028,
    "malicious_count": 65,
    "popular_threat_category": [
      {
        "value": "trojan",
        "count": 16
      },
      {
        "value": "worm",
        "count": 15
      }
    ],
    "popular_threat_name": [
      {
        "value": "fasong",
        "count": 28
      },
      {
        "value": "mint",
        "count": 6
      },
      {
        "value": "zard",
        "count": 6
      }
    ],
    "suggested_label": "trojan.fasong/mint",
    "av_family_vote": "win32",
    "has_family_label": true,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via polling"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1115",
        "signature_description": "read clipboard data"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1543",
        "signature_description": "create service"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1569",
        "signature_description": "create service"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "generate random numbers using the Delphi LCG",
      "create or open registry key",
      "create thread",
      "create or open file",
      "set file attributes",
      "delay execution",
      "PEB access",
      "resume thread",
      "encode data using XOR",
      "access PEB ldr_data",
      "allocate RWX memory",
      "hide graphical window",
      "set global application hook",
      "inspect load icon resource",
      "create TCP socket",
      "check mutex",
      "open process",
      "create process on Windows",
      "contain loop",
      "accept command line arguments",
      "query or enumerate registry value",
      "link function at runtime on Windows",
      "get geographical location",
      "read file on Windows",
      "write file on Windows",
      "get file size",
      "set thread local storage value",
      "get thread local storage value",
      "find graphical window",
      "delete file",
      "get disk size",
      "get OS version",
      "extract resource via kernel32 functions",
      "get system information on Windows",
      "read clipboard data",
      "capture screenshot",
      "log keystrokes via polling",
      "link many functions at runtime",
      "log keystrokes",
      "get keyboard layout",
      "set application hook",
      "enumerate gui resources",
      "get graphical window text",
      "encode data using Base64",
      "set registry value",
      "access the Windows event log",
      "get service handle",
      "create service",
      "delete service",
      "send data on socket",
      "send data",
      "receive data on socket",
      "receive data",
      "persist via Windows service",
      "resolve DNS",
      "get hostname",
      "terminate process",
      "get disk information",
      "create mutex",
      "enumerate files on Windows",
      "contain a resource (.rsrc) section",
      "packed with Mpress",
      "run as service",
      "reference Base64 string",
      "compiled with Borland Delphi",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\JVXUB.EXE"
      ],
      "registry_all_count": 7,
      "files_dropped": [
        "C:\\$Recycle.Bin\\ECR.EXE",
        "C:\\PerfLogs\\FXLM.EXE",
        "C:\\PerfLogs\\NFQQ.EXE",
        "C:\\Program Files (x86)\\EDZES.EXE",
        "C:\\Program Files (x86)\\JDHDD.EXE",
        "C:\\Program Files\\NCBE.EXE",
        "C:\\Program Files\\UQMH.EXE",
        "HJUV.EXE",
        "LFQRGF.EXE",
        "TPI.EXE",
        "WJF.EXE",
        "YHDP.EXE",
        "ZZB.EXE"
      ],
      "mutexes": [
        "ShimCacheMutex",
        "qqpass7"
      ],
      "processes": [
        "C:\\cuckoo\\ZZB.EXE"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}