# dataset/preprocessed/173ef63fcffc4f8a8e503ff9e66b09b822b1c6cf9b02d850df5e0b0314f940b5.json
{
  "hash": "173ef63fcffc4f8a8e503ff9e66b09b822b1c6cf9b02d850df5e0b0314f940b5",
  "report_text": "--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (13):\n  - encode data using XOR\n  - accept command line arguments\n  - get file size\n  - get disk size\n  - link function at runtime on Windows\n  - query or enumerate registry key\n  - query or enumerate registry value\n  - get common file path\n  - stop service\n  - get session user name\n  - get memory capacity\n  - check OS version\n  - packed with UPX\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - File System::Create File [C0016]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - Process::Create Thread [C0038]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Process [C0017]\n  - contain loop\n  - Execution::Command and Scripting Interpreter [E1059]\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - Discovery::File and Directory Discovery [E1083]\n  - Process::Set Thread Local Storage Value [C0041]\n  - get thread local storage value\n  - Process::Create Mutex [C0042]\n  - Data::Checksum::CRC32 [C0032.001]\n  - resolve function by parsing PE exports\n  - File System::Delete File [C0047]\n  - File System::Move File [C0063]\n  - Discovery::System Information Discovery [E1082]\n  - link function at runtime on Windows\n  - get OS version\n\n--- ARTIFACTS ---\n\nFILES WRITTEN (8):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\~rgeffere.tmp\n  C:\\Documents and Settings\\Administrator\\\u300c\u5f00\u59cb\u300d\u83dc\u5355\\\u7a0b\u5e8f\\\u542f\u52a8\\OAw9h0.exe\n  C:\\Documents and Settings\\All Users\\\u300c\u5f00\u59cb\u300d\u83dc\u5355\\\u7a0b\u5e8f\\\u542f\u52a8\\OAw9h0.exeDRAwD0vim1T.exe\n  C:\\Program Files\\Adobe\\Reader 9.0\\Reader\\AcroRd32.exe\n  C:\\Program Files\\Adobe\\Reader 9.0\\Reader\\AcroRd32Info.exe\n  C:\\Program Files\\WinRAR\\WinRAR.exe\n  C:\\WINDOWS\\6cmC7y.exe\n  C:\\WINDOWS\\system32\\shutdown.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\NPSSVC = C:\\WINDOWS\\6cmC7y.exe\u0000\n\nMUTEXES CREATED (4):\n  GhostBOT\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: CreateProcessA, GetProcAddress\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Packer] UPX\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n\nPE SECTIONS (entropy):\n  UPX0         entropy=5.76\n  UPX1         entropy=2.52\n  .rsrc        entropy=2.72\n  [overlay]    entropy=5.3597731590271  size=350160\n",
  "ground_truth": {
    "hash": "173ef63fcffc4f8a8e503ff9e66b09b822b1c6cf9b02d850df5e0b0314f940b5",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9275,
    "malicious_count": 64,
    "popular_threat_category": [
      {
        "value": "trojan",
        "count": 32
      },
      {
        "value": "worm",
        "count": 5
      }
    ],
    "popular_threat_name": [
      {
        "value": "gobot",
        "count": 14
      },
      {
        "value": "gbot",
        "count": 4
      },
      {
        "value": "rbot",
        "count": 2
      }
    ],
    "suggested_label": "trojan.gobot/gbot",
    "av_family_vote": "win32",
    "has_family_label": true,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry key"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1543",
        "signature_description": "stop service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open file",
      "generate random numbers using the Delphi LCG",
      "create thread",
      "encode data using XOR",
      "create or open registry key",
      "delay execution",
      "create process on Windows",
      "contain loop",
      "accept command line arguments",
      "read file on Windows",
      "write file on Windows",
      "get file size",
      "set thread local storage value",
      "get thread local storage value",
      "create mutex",
      "hash data with CRC32",
      "resolve function by parsing PE exports",
      "generate random numbers using a Mersenne Twister",
      "delete file",
      "move file",
      "get disk size",
      "link function at runtime on Windows",
      "get OS version",
      "query or enumerate registry key",
      "query or enumerate registry value",
      "set registry value",
      "get common file path",
      "copy file",
      "get service handle",
      "stop service",
      "execute anti-debugging instructions",
      "get session user name",
      "download URL",
      "receive data",
      "get memory capacity",
      "check OS version",
      "packed with UPX",
      "contain a resource (.rsrc) section",
      "contain an embedded PE file",
      "contains PDB path",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\NPSSVC"
      ],
      "registry_all_count": 1,
      "files_dropped": [],
      "mutexes": [
        "GhostBOT",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO"
  ]
}