# dataset/preprocessed/17e2085aa6f1516bf38cb75b0e015e9eb29fe21b75e81303310505763ab47875.json
{
  "hash": "17e2085aa6f1516bf38cb75b0e015e9eb29fe21b75e81303310505763ab47875",
  "report_text": "--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (19):\n  - contain obfuscated stackstrings\n  - check Internet connectivity via WinINet\n  - access PEB ldr_data\n  - hide graphical window\n  - encode data using XOR\n  - check if file exists\n  - get common file path\n  - accept command line arguments\n  - link function at runtime on Windows\n  - query or enumerate registry key\n  - delete registry value\n  - get hostname\n  - get file size\n  - get number of processors\n  - query or enumerate registry value\n  - persist via Windows service\n  - encrypt data using RC4 PRGA\n  - parse PE header\n  - link many functions at runtime\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - File System::Get File Attributes [C0049]\n  - File System::Create File [C0016]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - check Internet connectivity via WinINet\n  - access PEB ldr_data\n  - Process::Create Process [C0017]\n  - hide graphical window\n  - Process::Create Thread [C0038]\n  - Process::Resume Thread [C0054]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - contain loop\n  - resolve function by parsing PE exports\n  - Discovery::File and Directory Discovery [E1083]\n  - set current directory\n  - File System::Read File [C0051]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - link function at runtime on Windows\n  - Operating System::Registry::Set Registry Key [C0036.001]\n\nHIGHLIGHTED CALLS:\n  SetFileTime\n\n--- ARTIFACTS ---\n\nFILES DROPPED (1):\n  C:\\vglcon.pif  [sha256: ee9d4d56d8cff759...]\n\nPERSISTENCE REGISTRY KEYS (3):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\DisableNotifications = 1\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\DoNotAllowExceptions = 0\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\EnableFirewall = 0\n\nOTHER REGISTRY KEYS SET (10):\n  HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Explorer\\Advanced\\Hidden\n  HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\policies\\system\\DisableTaskMgr\n  HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\policies\\system\\DisableRegistryTools\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\AntiVirusOverride\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\AntiVirusDisableNotify\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\FirewallDisableNotify\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\FirewallOverride\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\UpdatesDisableNotify\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\UacDisableNotify\n  HKLM\\SOFTWARE\\Microsoft\\Security Center\\Svc\\AntiVirusOverride\n\nMUTEXES CREATED (39):\n  audiodg.exeM_932_\n  bssndrpt.exe.exeM_2420_\n  conhost.exeM_2960_\n  csrss.exeM_340_\n  csrss.exeM_396_\n  dwm.exeM_1248_\n  explorer.exeM_1268_\n  frida-winjector-helper-32.exeM_2928_\n  frida-winjector-helper-32.exeM_2972_\n  frida-winjector-helper-64.exeM_912_\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, LoadLibraryExW, LoadLibraryW, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  [Sign tool] Windows Authenticode\n\nPE SECTIONS (entropy):\n  .text        entropy=6.68\n  .rdata       entropy=4.84\n  .data        entropy=4.2\n  .rsrc        entropy=7.84\n  [overlay]    entropy=7.341536045074463  size=7184\n",
  "ground_truth": {
    "hash": "17e2085aa6f1516bf38cb75b0e015e9eb29fe21b75e81303310505763ab47875",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8451,
    "malicious_count": 60,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": null,
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry key"
      },
      {
        "id": "T1016",
        "signature_description": "check Internet connectivity via WinINet"
      },
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get hostname"
      },
      {
        "id": "T1083",
        "signature_description": "check if file exists"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry value"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1543",
        "signature_description": "persist via Windows service"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1569",
        "signature_description": "persist via Windows service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "PEB access",
      "get file attributes",
      "create or open file",
      "contain obfuscated stackstrings",
      "create or open registry key",
      "delay execution",
      "check Internet connectivity via WinINet",
      "access PEB ldr_data",
      "create process on Windows",
      "hide graphical window",
      "create thread",
      "resume thread",
      "encode data using XOR",
      "contain loop",
      "resolve function by parsing PE exports",
      "check if file exists",
      "set current directory",
      "get common file path",
      "read .ini file",
      "accept command line arguments",
      "link function at runtime on Windows",
      "set registry value",
      "query or enumerate registry key",
      "delete registry value",
      "extract resource via kernel32 functions",
      "write file on Windows",
      "print debug messages",
      "access the Windows event log",
      "get hostname",
      "resolve DNS",
      "initialize Winsock library",
      "get graphical window text",
      "create HTTP request",
      "connect to HTTP server",
      "send HTTP request",
      "send data",
      "read data from Internet",
      "receive data",
      "check HTTP status code",
      "get file size",
      "read file via mapping",
      "read file on Windows",
      "initialize hashing via WinCrypt",
      "hash data with MD5",
      "hash data via WinCrypt",
      "get number of processors",
      "query or enumerate registry value",
      "persist via Windows service",
      "encrypt data using RC4 PRGA",
      "enumerate PE sections",
      "hash data with CRC32",
      "parse PE header",
      "delete file",
      "get thread local storage value",
      "link many functions at runtime",
      "linked against XZip",
      "contain a resource (.rsrc) section",
      "linked against CPP standard library",
      "contains PDB path"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\DisableNotifications",
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\DoNotAllowExceptions",
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\EnableFirewall"
      ],
      "registry_all_count": 30,
      "files_dropped": [
        "C:\\vglcon.pif"
      ],
      "mutexes": [
        "audiodg.exeM_932_",
        "bssndrpt.exe.exeM_2420_",
        "conhost.exeM_2960_",
        "csrss.exeM_340_",
        "csrss.exeM_396_",
        "dwm.exeM_1248_",
        "explorer.exeM_1268_",
        "frida-winjector-helper-32.exeM_2928_",
        "frida-winjector-helper-32.exeM_2972_",
        "frida-winjector-helper-64.exeM_912_",
        "lsass.exeM_492_",
        "lsm.exeM_500_",
        "pythonw.exeM_1648_",
        "pyw.exeM_2316_",
        "searchfilterhost.exeM_2840_",
        "searchindexer.exeM_1460_",
        "searchprotocolhost.exeM_1356_",
        "services.exeM_484_",
        "smss.exeM_256_",
        "spoolsv.exeM_1056_",
        "sppsvc.exeM_2216_",
        "svchost.exeM_108_",
        "svchost.exeM_1100_",
        "svchost.exeM_1288_",
        "svchost.exeM_1484_",
        "svchost.exeM_2404_",
        "svchost.exeM_300_",
        "svchost.exeM_612_",
        "svchost.exeM_688_",
        "svchost.exeM_756_",
        "svchost.exeM_828_",
        "svchost.exeM_864_",
        "taskhost.exeM_1176_",
        "taskhost.exeM_1604_",
        "uxJLpe1m",
        "wininit.exeM_380_",
        "winlogon.exeM_436_",
        "wmiprvse.exeM_2640_",
        "wmpnetwk.exeM_644_"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Jujubox"
  ]
}