# dataset/preprocessed/18a98d126c90020966671bf6a8435fedbc2493ba98ba76b03048df3a845267e0.json
{
  "hash": "18a98d126c90020966671bf6a8435fedbc2493ba98ba76b03048df3a845267e0",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (3):\n  ((null)) C:\\18a98d126c90020966671bf6a8435fedbc2493ba98ba76b03048df3a845267e0(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\18a98d126c900209_russian nude blowjob full movie feet .mpg.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american animal lesbian catfight ash .mpg.exe  [sha256: 8fe8bb3e5cb6ab1d...]\n  american cum lesbian masturbation glans stockings .avi.exe  [sha256: aceb107a54a3d3d4...]\n  beast [bangbus] cock hotel .rar.exe  [sha256: 4dac5e8d05ec9420...]\n  brasilian cum xxx uncut granny .rar.exe  [sha256: def30c690277e25d...]\n  brasilian kicking lesbian [milf]  (Tatjana).mpeg.exe  [sha256: 61f9596a669e8887...]\n  canadian horse licking glans penetration  (Sylvia).mpeg.exe  [sha256: 1d481dd490a2084c...]\n  danish horse beast big glans sm  (Melissa).avi.exe  [sha256: 2ad40a07992cba1d...]\n  deinstall.exe  [sha256: c290ed56f1913ff4...]\n  driver.exe  [sha256: b38b2d47f10229bd...]\n  drwatson32.exe  [sha256: f49e412d8829f299...]\n  freeftp.exe  [sha256: abaef2d28f1860cf...]\n  fucking [milf] hole .mpeg.exe  [sha256: 987887288858ed3f...]\n  gay lesbian young .mpeg.exe  [sha256: 7bf0270cfeaab342...]\n  IE9setup.exe  [sha256: aa0898f8316cb31f...]\n  indian kicking fucking [free] titts leather .mpg.exe  [sha256: faeb2d0e92cf881e...]\n  indian porn lesbian [free]  (Janette).zip.exe  [sha256: 0332d31f35b4d0fb...]\n  indian porn trambling public  (Sylvia).zip.exe  [sha256: 7650f3c9707dd2ea...]\n  italian animal trambling [milf] .avi.exe  [sha256: b635a6b6128c0bd7...]\n  italian cumshot lesbian catfight gorgeoushorny .mpeg.exe  [sha256: aba3f0268ed01449...]\n  lingerie catfight leather .mpeg.exe  [sha256: 0b7f47dc49e17c70...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\indian porn lesbian [free]  (Janette).zip.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\russian gang bang xxx uncut feet 50+ .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\danish horse beast big glans sm  (Melissa).avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\tyrkish kicking hardcore big feet .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\swedish porn lingerie [free] hole .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\italian cumshot lesbian catfight gorgeoushorny .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\gay lesbian young .mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\beast [bangbus] cock hotel .rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\russian nude trambling [free] hole hotel .rar.exe\n  C:\\Program Files\\Movie Maker\\Shared\\trambling hot (!) beautyfull .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\brasilian kicking lesbian [milf]  (Tatjana).mpeg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\tyrkish kicking trambling several models glans .mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\russian action horse hot (!) .mpg.exe\n  C:\\WINDOWS\\Temp\\brasilian cum xxx uncut granny .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\indian porn trambling public  (Sylvia).zip.exe\n  C:\\WINDOWS\\assembly\\temp\\canadian horse licking glans penetration  (Sylvia).mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\lingerie public cock 40+  (Melissa).zip.exe\n  C:\\WINDOWS\\ime\\shared\\american cum lesbian masturbation glans stockings .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\security\\templates\\norwegian sperm hot (!) young .rar.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880061149597168  size=362124\n",
  "ground_truth": {
    "hash": "18a98d126c90020966671bf6a8435fedbc2493ba98ba76b03048df3a845267e0",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9306,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\swedish gang bang sperm [free] feet wifey  (Karin).mpg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\blowjob hot (!) hole  (Sonja,Jade).avi.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\hardcore sleeping .avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\swedish beastiality trambling licking  (Karin).avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\lesbian sleeping titts 40+ .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\tyrkish kicking lesbian licking .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\american cum fucking [bangbus] leather .zip.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\italian cumshot beast lesbian  (Jade).zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\tyrkish handjob fucking [free] .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\japanese kicking gay lesbian redhair .zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\fucking hot (!) hole granny .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\beast masturbation glans .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\black beastiality trambling [milf] glans balls .mpeg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american animal lesbian catfight ash .mpg.exe",
        "american cum lesbian masturbation glans stockings .avi.exe",
        "beast [bangbus] cock hotel .rar.exe",
        "brasilian cum xxx uncut granny .rar.exe",
        "brasilian kicking lesbian [milf]  (Tatjana).mpeg.exe",
        "canadian horse licking glans penetration  (Sylvia).mpeg.exe",
        "danish horse beast big glans sm  (Melissa).avi.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking [milf] hole .mpeg.exe",
        "gay lesbian young .mpeg.exe",
        "indian kicking fucking [free] titts leather .mpg.exe",
        "indian porn lesbian [free]  (Janette).zip.exe",
        "indian porn trambling public  (Sylvia).zip.exe",
        "italian animal trambling [milf] .avi.exe",
        "italian cumshot lesbian catfight gorgeoushorny .mpeg.exe",
        "lingerie catfight leather .mpeg.exe",
        "lingerie public cock 40+  (Melissa).zip.exe",
        "mssrv.exe",
        "norwegian sperm hot (!) young .rar.exe",
        "russian action horse hot (!) .mpg.exe",
        "russian gang bang xxx uncut feet 50+ .zip.exe",
        "russian kicking lesbian lesbian gorgeoushorny .zip.exe",
        "russian nude trambling [free] hole hotel .rar.exe",
        "setup.exe",
        "soccer.exe",
        "sperm [milf] feet high heels .rar.exe",
        "start.exe",
        "swedish fetish beast full movie .mpeg.exe",
        "swedish porn lingerie [free] hole .mpg.exe",
        "trambling hot (!) beautyfull .zip.exe",
        "tyrkish kicking hardcore big feet .zip.exe",
        "tyrkish kicking trambling several models glans .mpeg.exe",
        "undelete.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}