# dataset/preprocessed/1a329d4a2699a0d4be7a482d4261315c81f4ab81dfc86b4a4382b734982c7e18.json
{
  "hash": "1a329d4a2699a0d4be7a482d4261315c81f4ab81dfc86b4a4382b734982c7e18",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  C:\\Users\\Lucas\\AppData\\Local\\Temp\\wviZi.exe\n  C:\\Windows\\SysWOW64\\cmd.exe\n  C:\\Windows\\SysWOW64\\cmd.exe\n  C:\\Windows\\SysWOW64\\reg.exe\n  C:\\Windows\\SysWOW64\\reg.exe\n  C:\\windows\\temp\\10022.exe\n\nCOMMAND EXECUTIONS (9):\n  C:\\Users\\Lucas\\AppData\\Local\\Temp\\wviZi.exe\n  cmd /c reg ADD HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v freizer /t REG_SZ /d C:\\WINDOWS\\system32\\freizer.exe /f\n  cmd /c reg ADD HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v svchost /t REG_SZ /d C:\\WINDOWS\\system32\\svchost.exe /f\n  reg ADD HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v freizer /t REG_SZ /d C:\\WINDOWS\\system32\\freizer.exe /f\n  reg ADD HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v svchost /t REG_SZ /d C:\\WINDOWS\\system32\\svchost.exe /f\n  C:\\windows\\temp\\10022.exe\n  \"C:\\windows\\temp\\10022.exe\"\n  reg  ADD HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v freizer /t REG_SZ /d C:\\WINDOWS\\system32\\freizer.exe /f\n  reg  ADD HKEY_LOCAL_MACHINE\\Software\\Microsoft\\Windows\\CurrentVersion\\Run /v svchost /t REG_SZ /d C:\\WINDOWS\\system32\\svchost.exe /f\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (9):\n  - accept command line arguments\n  - get file size\n  - query or enumerate registry value\n  - get geographical location\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - get common file path\n  - persist via Run registry key\n  - reference Base64 string\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - File System::Create File [C0016]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - reference absolute stream path on Windows\n  - Process::Create Process [C0017]\n  - contain loop\n  - Execution::Command and Scripting Interpreter [E1059]\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - Discovery::File and Directory Discovery [E1083]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - get geographical location\n  - get OS version\n  - Process::Set Thread Local Storage Value [C0041]\n  - get thread local storage value\n  - Process::Create Mutex [C0042]\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - File System::Delete File [C0047]\n  - File System::Copy File [C0045]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (3):\n  \\Windows\\Temp\\10022.exe  [sha256: e89ac835ed17dbf4...]\n  \\Windows\\svchost.exe  [sha256: 1a329d4a2699a0d4...]\n  \\Program Files\\Adobe\\Acrobat Reader DC\\Reader\\ADelRCP.exe  [sha256: f187da932a68d8d5...]\n\nFILES WRITTEN (4):\n  C:\\WINDOWS\\svchost.exe\n  C:\\WINDOWS\\system32\\freizer.exe\n  C:\\Windows\\system32\\smrss.exe\n  C:\\tmpmqgko7\\44344.exe\n\nPERSISTENCE REGISTRY KEYS (5):\n  HKLM\\SOFTWARE\\MICROSOFT\\WINDOWS NT\\CURRENTVERSION\\WINLOGON = shell\n  HKLM\\SOFTWARE\\WOW6432NODE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\RUN = freizer\n  HKLM\\SOFTWARE\\WOW6432NODE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\RUN = svchost\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\freizer = C\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\svchost = C\n\nOTHER REGISTRY KEYS SET (3):\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\ZONEMAP\n  \\REGISTRY\\MACHINE\\HexStream = 53006F006600740077006100720065005C004D006900630072006F0073006F00660074005C00570069006E0064\n  \\REGISTRY\\USER\\S-1-5-21-3616907901-860437140-2684471051-1001\\HexStream = 53006F006600740077006100720065005C004D006900630\n\nMUTEXES CREATED (5):\n  BELARUS-VIRUS-MAKER\n  Local\\ZoneAttributeCacheCounterMutex\n  Local\\ZonesCacheCounterMutex\n  Local\\ZonesCounterMutex\n  Local\\ZonesLockedCacheCounterMutex\n\nMUTEXES OPENED (1):\n  Local\\MSCTF.Asm.MutexDefault1\n\nPROCESSES CREATED (7):\n  10022.exe\n  C:\\Users\\Lucas\\AppData\\Local\\Temp\\wviZi.exe\n  C:\\Windows\\SysWOW64\\cmd.exe\n  C:\\Windows\\SysWOW64\\reg.exe\n  C:\\windows\\temp\\10022.exe\n  cmd.exe\n  reg.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  kernel32.dll: GetProcAddress, VirtualAlloc, WinExec\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [sfx] temporary EXE SFX\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  [Overlay] EXE file\n  Packers: {'PEiD': 'BobSoft Mini Delphi -> BoB / BobSoft'}\n\nPE SECTIONS (entropy):\n  CODE         entropy=6.47\n  DATA         entropy=1.95\n  BSS          entropy=0.0\n  .idata       entropy=3.71\n  .tls         entropy=0.0\n  .rdata       entropy=0.2\n  .reloc       entropy=5.68\n  .rsrc        entropy=2.88\n  [overlay]    entropy=1.8560158014297485  size=16692478\n",
  "ground_truth": {
    "hash": "1a329d4a2699a0d4be7a482d4261315c81f4ab81dfc86b4a4382b734982c7e18",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9014,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "malware",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "reference Base64 string"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "generate random numbers using the Delphi LCG",
      "create or open file",
      "create or open registry key",
      "reference absolute stream path on Windows",
      "create process on Windows",
      "contain loop",
      "accept command line arguments",
      "read file on Windows",
      "write file on Windows",
      "get file size",
      "query or enumerate registry value",
      "get geographical location",
      "get OS version",
      "set thread local storage value",
      "get thread local storage value",
      "create mutex",
      "link function at runtime on Windows",
      "link many functions at runtime",
      "delete file",
      "copy file",
      "get common file path",
      "set registry value",
      "persist via Run registry key",
      "contain an embedded PE file",
      "contain a thread local storage (.tls) section",
      "contain a resource (.rsrc) section",
      "reference Base64 string",
      "packaged as a NSIS installer",
      "contains PDB path",
      "compiled with Borland Delphi",
      "(internal) installer file limitation",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\MICROSOFT\\WINDOWS NT\\CURRENTVERSION\\WINLOGON",
        "HKLM\\SOFTWARE\\WOW6432NODE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\RUN",
        "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\freizer",
        "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\svchost"
      ],
      "registry_all_count": 7,
      "files_dropped": [
        "\\Program Files\\Adobe\\Acrobat Reader DC\\Reader\\ADelRCP.exe",
        "\\Windows\\Temp\\10022.exe",
        "\\Windows\\svchost.exe"
      ],
      "mutexes": [
        "BELARUS-VIRUS-MAKER",
        "Local\\ZoneAttributeCacheCounterMutex",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "10022.exe",
        "C:\\Users\\Lucas\\AppData\\Local\\Temp\\wviZi.exe",
        "C:\\Windows\\SysWOW64\\cmd.exe",
        "C:\\Windows\\SysWOW64\\reg.exe",
        "C:\\windows\\temp\\10022.exe",
        "cmd.exe",
        "reg.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Lastline",
    "SNDBOX",
    "VirusTotal Jujubox"
  ]
}