# dataset/preprocessed/1aef8163f2bf3f2c5c116c068ea9f5780820b606db396f55dc0c6cae3cff0854.json
{
  "hash": "1aef8163f2bf3f2c5c116c068ea9f5780820b606db396f55dc0c6cae3cff0854",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nIP TRAFFIC (1):\n  <MACHINE_DNS_SERVER>:53 (UDP)\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  1aef8163f2bf3f2c_horse several models glans .rar.exe\n    C:\\Users\\<USER>\\Downloads\\1aef8163f2bf3f2c_horse several models glans .rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\1aef8163f2bf3f2c5c116c068ea9f5780820b606db396f55dc0c6cae3cff0854(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\1aef8163f2bf3f2c_horse several models glans .rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american cumshot bukkake [bangbus] black hairunshaved .mpeg.exe  [sha256: 8e3df2a52de3b742...]\n  american fetish beast lesbian  (Tatjana).mpeg.exe  [sha256: 1c0cd69648f45017...]\n  blowjob several models .avi.exe  [sha256: b318baa3165034a6...]\n  brasilian beastiality bukkake hot (!) ash  (Anniston,Karin).avi.exe  [sha256: 46fc2767a5500711...]\n  brasilian horse fucking catfight feet .mpeg.exe  [sha256: 4126b4a89b1e0c49...]\n  bukkake catfight ejaculation .mpeg.exe  [sha256: ebe94c75783729d6...]\n  bukkake several models .avi.exe  [sha256: e09125c123c43b8d...]\n  danish action horse [milf] feet high heels  (Sylvia).rar.exe  [sha256: a43fcafc1e422f59...]\n  defrag.exe  [sha256: f3755c9c4705fdf2...]\n  defrag64.exe  [sha256: 562bc4be5d34923f...]\n  driver.exe  [sha256: 3976cc5059440b25...]\n  drwatson32.exe  [sha256: 1e1610832b0c8e59...]\n  freeftp.exe  [sha256: f0d1147a4800a09d...]\n  fucking full movie bondage .mpg.exe  [sha256: 1fcba9d5535bfd51...]\n  gay licking  (Curtney).zip.exe  [sha256: 5f4b3a1c8b40782a...]\n  gay [free] .zip.exe  [sha256: dcc05ae4cc7fdde9...]\n  hardcore several models beautyfull .avi.exe  [sha256: 520385e70d5b3450...]\n  indian beastiality hardcore catfight titts .zip.exe  [sha256: e88c2dfcfc006d0b...]\n  indian fetish xxx hidden cock .mpeg.exe  [sha256: cd7ca69425001db6...]\n  indian gang bang trambling several models glans granny .rar.exe  [sha256: 0df71a7b39f9a808...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\brasilian beastiality bukkake hot (!) ash  (Anniston,Karin).avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\lingerie sleeping  (Sarah).avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\russian handjob trambling lesbian granny .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\danish action horse [milf] feet high heels  (Sylvia).rar.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\indian fetish xxx hidden cock .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\tyrkish cum sperm [milf] glans  (Jenna,Janette).mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\trambling [free]  (Curtney).rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\blowjob several models .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\xxx several models hole redhair .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\xxx hot (!) cock leather  (Sylvia).avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\indian beastiality hardcore catfight titts .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\american cumshot bukkake [bangbus] black hairunshaved .mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\bukkake catfight ejaculation .mpeg.exe\n  C:\\WINDOWS\\Temp\\gay [free] .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\gay licking  (Curtney).zip.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\trambling full movie latex .mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\brasilian horse fucking catfight feet .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\italian gang bang hardcore catfight pregnant .zip.exe\n  C:\\WINDOWS\\ime\\shared\\lesbian girls hole .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880110263824463  size=707960\n",
  "ground_truth": {
    "hash": "1aef8163f2bf3f2c5c116c068ea9f5780820b606db396f55dc0c6cae3cff0854",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [
        "<MACHINE_DNS_SERVER>"
      ],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\action nude [bangbus] .mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\cumshot sleeping fishy  (Gina,Britney).zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\beast sleeping balls  (Jenna).mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\french horse beast voyeur feet .avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\indian horse masturbation gorgeoushorny .zip.exe",
        "C:\\Program Files\\dotnet\\shared\\chinese beast bukkake several models fishy .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "american cumshot bukkake [bangbus] black hairunshaved .mpeg.exe",
        "american fetish beast lesbian  (Tatjana).mpeg.exe",
        "blowjob several models .avi.exe",
        "brasilian beastiality bukkake hot (!) ash  (Anniston,Karin).avi.exe",
        "brasilian horse fucking catfight feet .mpeg.exe",
        "bukkake catfight ejaculation .mpeg.exe",
        "bukkake several models .avi.exe",
        "danish action horse [milf] feet high heels  (Sylvia).rar.exe",
        "defrag.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking full movie bondage .mpg.exe",
        "gay [free] .zip.exe",
        "gay licking  (Curtney).zip.exe",
        "hardcore several models beautyfull .avi.exe",
        "indian beastiality hardcore catfight titts .zip.exe",
        "indian fetish xxx hidden cock .mpeg.exe",
        "indian gang bang trambling several models glans granny .rar.exe",
        "italian gang bang hardcore catfight pregnant .zip.exe",
        "japanese fetish blowjob [bangbus] lady .avi.exe",
        "lesbian girls hole .rar.exe",
        "lingerie sleeping  (Sarah).avi.exe",
        "mssrv.exe",
        "readme.txt.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}