# dataset/preprocessed/1c0bcd12710a2e80a1d9773d04669de15421660b57f722343193287b0e34aadb.json
{
  "hash": "1c0bcd12710a2e80a1d9773d04669de15421660b57f722343193287b0e34aadb",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  cumshot lesbian ash .zip.exe\n    C:\\Users\\<USER>\\Downloads\\cumshot lesbian ash .zip.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\1c0bcd12710a2e80a1d9773d04669de15421660b57f722343193287b0e34aadb(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\cumshot lesbian ash .zip.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  action [free] .mpg.exe  [sha256: c9048e14bece61a2...]\n  asian kicking lesbian .rar.exe  [sha256: e3999b01803a1392...]\n  beast girls feet black hairunshaved .avi.exe  [sha256: 6327b35b73845215...]\n  blowjob lesbian young .mpg.exe  [sha256: 8330151d4487081f...]\n  bukkake blowjob uncut granny .rar.exe  [sha256: df8d182b975ddff1...]\n  bukkake trambling full movie latex  (Ashley,Britney).avi.exe  [sha256: a87527f23f785302...]\n  bukkake voyeur .mpg.exe  [sha256: e17f51207ea05bcd...]\n  canadian blowjob xxx [milf] girly .zip.exe  [sha256: 32254c3e84a93669...]\n  canadian handjob masturbation vagina balls .mpeg.exe  [sha256: bf478a9dd225de40...]\n  canadian xxx voyeur 40+ .rar.exe  [sha256: b4023c641db7a3af...]\n  cum masturbation  (Sylvia).mpeg.exe  [sha256: 290b6942790fa02b...]\n  danish gay big fishy .mpeg.exe  [sha256: b833b17139a18a58...]\n  defrag.exe  [sha256: 3eaf972ecc2eef67...]\n  defrag64.exe  [sha256: e72c13aa3c794994...]\n  driver.exe  [sha256: d5ce98b336f52e3c...]\n  drwatson32.exe  [sha256: 33f03111ed9ea5a2...]\n  edit.exe  [sha256: 10a07aff2e223bc4...]\n  explorer.exe  [sha256: bd9327a285447596...]\n  freeftp.exe  [sha256: ec37fca5f093419e...]\n  french xxx kicking uncut hole .rar.exe  [sha256: df83cdfd1d741596...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\horse fucking public nipples .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\asian kicking lesbian .rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\danish gay big fishy .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\sperm [milf] .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\tyrkish beastiality horse [milf] .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\lesbian sleeping .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\cum masturbation  (Sylvia).mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\lingerie action several models ejaculation  (Gina).rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\kicking uncut gorgeoushorny .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\hardcore [free] pregnant .avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\action [free] .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\bukkake voyeur .mpg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\blowjob lesbian young .mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\spanish horse public  (Melissa).avi.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\bukkake blowjob uncut granny .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\sperm voyeur titts granny .mpg.exe\n  C:\\WINDOWS\\assembly\\tmp\\japanese nude horse big vagina .avi.exe\n  C:\\WINDOWS\\ime\\shared\\italian horse cum [milf] sm .mpg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\french xxx kicking uncut hole .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.64\n  [overlay]    entropy=7.880005359649658  size=429891\n",
  "ground_truth": {
    "hash": "1c0bcd12710a2e80a1d9773d04669de15421660b57f722343193287b0e34aadb",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9444,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\sperm girls  (Janette,Jenna).mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\spanish horse beast [free] ash ash .rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\sperm beast hot (!) .rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\french lingerie kicking voyeur leather  (Jenna).zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\beast beast full movie glans high heels  (Britney,Sylvia).mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\gang bang public 50+  (Melissa).avi.exe",
        "C:\\Windows\\mssrv.exe",
        "action [free] .mpg.exe",
        "asian kicking lesbian .rar.exe",
        "beast girls feet black hairunshaved .avi.exe",
        "blowjob lesbian young .mpg.exe",
        "bukkake blowjob uncut granny .rar.exe",
        "bukkake trambling full movie latex  (Ashley,Britney).avi.exe",
        "bukkake voyeur .mpg.exe",
        "canadian blowjob xxx [milf] girly .zip.exe",
        "canadian handjob masturbation vagina balls .mpeg.exe",
        "canadian xxx voyeur 40+ .rar.exe",
        "cum masturbation  (Sylvia).mpeg.exe",
        "danish gay big fishy .mpeg.exe",
        "defrag.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "edit.exe",
        "explorer.exe",
        "freeftp.exe",
        "french xxx kicking uncut hole .rar.exe",
        "hardcore [free] pregnant .avi.exe",
        "horse fucking public nipples .avi.exe",
        "indian kicking porn public legs .mpg.exe",
        "italian horse cum [milf] sm .mpg.exe",
        "japanese kicking fucking full movie sm .avi.exe",
        "japanese nude horse big vagina .avi.exe",
        "kicking uncut gorgeoushorny .avi.exe",
        "lesbian sleeping .mpeg.exe",
        "lingerie action several models ejaculation  (Gina).rar.exe",
        "malaysia trambling full movie upskirt .rar.exe",
        "mssrv.exe",
        "rundll64.exe",
        "russian animal several models upskirt .mpg.exe",
        "setup.exe",
        "spanish horse public  (Melissa).avi.exe",
        "sperm [milf] .zip.exe",
        "sperm voyeur titts granny .mpg.exe",
        "start.exe",
        "tyrkish beastiality horse [milf] .zip.exe",
        "uninstall.exe",
        "word.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}