# dataset/preprocessed/1c72045025e99eaf94a21de8b25750eb713412106b562174d753d26ba8ef3fd3.json
{
  "hash": "1c72045025e99eaf94a21de8b25750eb713412106b562174d753d26ba8ef3fd3",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (8):\n  ftp.7host.com\n  ftp.7host.com \u2192 \n  ftp.alizametal.com.tr \u2192 185.42.173.103\n  ftp.freebiri.t15.org \u2192 \n  ftp.solisiletisim.com \u2192 \n  ftp.yesimcopy.com \u2192 89.43.28.59\n  gripto.leadhoster.com\n  gripto.leadhoster.com \u2192 127.0.0.1\n\nIP TRAFFIC (1):\n  8.8.4.4:53 (UDP)\n\nSERVICES OPENED (1):\n  RASMAN\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  38a54bcf9f4559da7bfc5ebcfa6de130.exe\n    73b12ce12dea223e3f12e0e19b669c54.exe\n  d1c56374fff0243832b8696d133b7861.exe\n  d1c56374fff0243832b8696d133b7861.exe\n    7c34a78b4422b61608ab858fb3833bc4.exe\n    4b7257ed5fa55aad7c4290b0188dc872.exe\n  d1c56374fff0243832b8696d133b7861.exe\n\nCOMMAND EXECUTIONS (4):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n  C:\\Program Files (x86)\\a4685e47\\jusched.exe \n  C:\\Windows\\system32\\DllHost.exe /Processid:{F9717507-6651-4EDB-BFF7-AE615179BCCF}\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - enumerate files on Windows\n  - parse PE header\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - Discovery::File and Directory Discovery [E1083]\n  - parse PE header\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - contain a resource (.rsrc) section\n  - debug build\n\nHIGHLIGHTED CALLS:\n  C:\\Program Files (x86)\\a4685e47\\jusched.exe\n  C:\\Windows\\system32\\DllHost.exe /Processid:{F9717507-6651-4EDB-BFF7-AE615179BCCF}\n  kernel32.dll.FreeLibrary\n  kernel32.dll.FindFirstFileA\n  kernel32.dll.WriteFile\n  kernel32.dll.ReadFile\n  shell32.dll.ShellExecuteA\n  wininet.dll.InternetConnectA\n  samlib.dll.SamConnect\n  kernel32.dll.RegOpenKeyExW\n\n--- ARTIFACTS ---\n\nFILES DROPPED (2):\n  jusched.exe  [sha256: 4e0c58e2c471c04f...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: 18758f3868931dbc...]\n\nFILES WRITTEN (30):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files (x86)\\71d2c87e\\info_a\n  C:\\Program Files (x86)\\a4685e47\\a4685e47\n  C:\\Program Files (x86)\\a4685e47\\info_a\n  C:\\Program Files (x86)\\a4685e47\\jusched.exe\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\info_a\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\Program Files\\f440491e\\f440491e\n  C:\\Program Files\\f440491e\\info_a\n  C:\\Program Files\\f440491e\\jusched.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Microsoft\\Windows\\Temporary Internet Files\\counters.dat\n  C:\\WINDOWS\\Tasks\\2tdU3eap.job\n  C:\\Windows\\SoftwareDistribution\\DataStore\\DataStore.edb\n  C:\\Windows\\SoftwareDistribution\\DataStore\\Logs\\edb.chk\n  C:\\Windows\\Tasks\\2tdU3eap.job\n  C:\\Windows\\Temp\\fwtsqmfile03.sqm\n  C:\\Windows\\appcompat\\Programs\\RecentFileCache.bcf\n  C:\\Windows\\sysnative\\Tasks\\2tdU3eap\n  C:\\Windows\\sysnative\\wbem\\repository\\INDEX.BTR\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:edaJ2aUa3tpv\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe = C:\\Program Files\\f440491e\\jusched.exe:*:Enabled:edaJ2aUa3tpv\u0000\n\nOTHER REGISTRY KEYS SET (2):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Program \n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\n\nMUTEXES CREATED (11):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  RasPbFile\n  ShimCacheMutex\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n\nMUTEXES OPENED (4):\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\nPROCESSES CREATED (1):\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.49\n  .rdata       entropy=4.64\n  .data        entropy=0.99\n  .idata       entropy=3.28\n  .rsrc        entropy=1.81\n  .reloc       entropy=0.0\n  [overlay]    entropy=6.6992058753967285  size=6982\n",
  "ground_truth": {
    "hash": "1c72045025e99eaf94a21de8b25750eb713412106b562174d753d26ba8ef3fd3",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9014,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "enumerate files on Windows",
      "parse PE header",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "contain a resource (.rsrc) section",
      "debug build"
    ],
    "iocs": {
      "ips": [
        "8.8.4.4"
      ],
      "domains": [
        "ftp.7host.com",
        "ftp.alizametal.com.tr",
        "ftp.freebiri.t15.org",
        "ftp.solisiletisim.com",
        "ftp.yesimcopy.com",
        "gripto.leadhoster.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe"
      ],
      "registry_all_count": 4,
      "files_dropped": [
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "jusched.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "RasPbFile",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \""
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox",
    "Yomi Hunter"
  ]
}