# dataset/preprocessed/1f04633e703e7840e844af657a94c1150855d74cbb066bd026a7612f12d9af8b.json
{
  "hash": "1f04633e703e7840e844af657a94c1150855d74cbb066bd026a7612f12d9af8b",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  ftp.byethost12.com \u2192 185.27.134.11\n\nIP TRAFFIC (1):\n  185.27.134.11:21 (TCP)\n\nSERVICES OPENED (1):\n  RASMAN\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  <PATH_SAMPLE.EXE>\n    %ProgramFiles(x86)%\\afddb486\\jusched.exe\n  <SYSTEM32>\\svchost.exe\n  %ProgramFiles(x86)%\\afddb486\\jusched.exe\n\nCOMMAND EXECUTIONS (2):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - parse PE header\n  - enumerate files on Windows\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - parse PE header\n  - Discovery::File and Directory Discovery [E1083]\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - debug build\n  - contain a resource (.rsrc) section\n\n--- ARTIFACTS ---\n\nFILES DROPPED (6):\n  jusched.exe  [sha256: 10fc85e914468c6f...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: eae4a940db353d1b...]\n  %ProgramFiles(x86)%\\afddb486\\afddb486  [sha256: 54f9150d1268f7b4...]\n  %ProgramFiles(x86)%\\afddb486\\jusched.exe  [sha256: e6b831a2864de0bb...]\n  <SYSTEM32>\\tasks\\update23  [sha256: a93e4bb1e7785005...]\n  %WINDIR%\\tasks\\update23.job  [sha256: e9e4ceadc8289f59...]\n\nFILES WRITTEN (5):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\WINDOWS\\Tasks\\Update23.job\n  C:\\Windows\\Tasks\\Update23.job\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:JavaUpdate23\n\nOTHER REGISTRY KEYS SET (10):\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\UNCAsIntranet\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\AutoDetect\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectId_\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\LruList\\CurrentLru\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectId\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectLru\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectLru_\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\ObjectTable\\E9\\_FileId_\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\ObjectTable\\E9\\_Usn_\n  \\REGISTRY\\A\\{b19ab72e-fa9a-11e9-aa78-5e61e6b2e535}\\DefaultObjectStore\\ObjectTable\\E9\\_UsnJournalId_\n\nMUTEXES CREATED (2):\n  RasPbFile\n  ShimCacheMutex\n\nPROCESSES CREATED (3):\n  %ProgramFiles(x86)%\\afddb486\\jusched.exe\n  <PATH_SAMPLE.EXE>\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.88\n  .rdata       entropy=4.35\n  .data        entropy=0.44\n  .idata       entropy=2.9\n  .rsrc        entropy=2.5\n  .reloc       entropy=0.5\n  [overlay]    entropy=6.916566371917725  size=10382\n",
  "ground_truth": {
    "hash": "1f04633e703e7840e844af657a94c1150855d74cbb066bd026a7612f12d9af8b",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9306,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "parse PE header",
      "enumerate files on Windows",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "debug build",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [
        "185.27.134.11"
      ],
      "domains": [
        "ftp.byethost12.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe"
      ],
      "registry_all_count": 34,
      "files_dropped": [
        "%ProgramFiles(x86)%\\afddb486\\afddb486",
        "%ProgramFiles(x86)%\\afddb486\\jusched.exe",
        "%WINDIR%\\tasks\\update23.job",
        "<SYSTEM32>\\tasks\\update23",
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "jusched.exe"
      ],
      "mutexes": [
        "RasPbFile",
        "ShimCacheMutex"
      ],
      "processes": [
        "%ProgramFiles(x86)%\\afddb486\\jusched.exe",
        "<PATH_SAMPLE.EXE>",
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \""
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Dr.Web vxCube",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}