# dataset/preprocessed/207b3e9e6c9f35e7c4b9a830d73cf975ab524ec75fdf1fb33d022e4f0dbc5275.json
{
  "hash": "207b3e9e6c9f35e7c4b9a830d73cf975ab524ec75fdf1fb33d022e4f0dbc5275",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (4):\n  ftp.0adz.com \u2192 45.56.79.23\n  ftp.0adz.com \u2192 96.126.123.244, 45.33.23.183\n  ftp.search.instantfreesite.com\n  gripto.leadhoster.com \u2192 127.0.0.1\n\nIP TRAFFIC (2):\n  45.33.23.183:21 (TCP)\n  96.126.123.244:21 (TCP)\n\nSERVICES OPENED (2):\n  RASMAN\n  Sens\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  jusched.exe\n    C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n\nCOMMAND EXECUTIONS (4):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n  \"C:\\Program Files\\48db4bad\\jusched.exe\" \n   \n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - enumerate files on Windows\n  - parse PE header\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - Discovery::File and Directory Discovery [E1083]\n  - parse PE header\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - debug build\n  - contain a resource (.rsrc) section\n\nHIGHLIGHTED CALLS:\n  SetFileTime\n\n--- ARTIFACTS ---\n\nFILES DROPPED (2):\n  jusched.exe  [sha256: e36f52948f4b421c...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: b5047e8cda63bc69...]\n\nFILES WRITTEN (11):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files (x86)\\71d2c87e\\info_a\n  C:\\Program Files\\48db4bad\n  C:\\Program Files\\48db4bad\\48db4bad\n  C:\\Program Files\\48db4bad\\info_a\n  C:\\Program Files\\48db4bad\\jusched.exe\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\info_a\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\WINDOWS\\Tasks\\2tdU3eap.job\n  C:\\Windows\\Tasks\\2tdU3eap.job\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:edaJ2aUa3tpv\n\nOTHER REGISTRY KEYS SET (4):\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\ProxyEnable\n  HKU\\S-1-5-21-364843204-231886559-199882026-1001\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\\n  HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\AutoDetect\n  HKEY_CURRENT_USER\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\UNCAsIntranet\n\nMUTEXES CREATED (2):\n  RasPbFile\n  ShimCacheMutex\n\nMUTEXES OPENED (3):\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!history!history.ie5!\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  Local\\c:!users!<USER>!appdata!roaming!microsoft!windows!cookies!\n\nPROCESSES CREATED (2):\n  C:\\Program Files\\48db4bad\\jusched.exe\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.66\n  .rdata       entropy=4.65\n  .data        entropy=0.99\n  .idata       entropy=3.27\n  .rsrc        entropy=1.81\n  .reloc       entropy=0.0\n  .rdata       entropy=0.0\n  [overlay]    entropy=6.935023784637451  size=2718\n",
  "ground_truth": {
    "hash": "207b3e9e6c9f35e7c4b9a830d73cf975ab524ec75fdf1fb33d022e4f0dbc5275",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9,
    "malicious_count": 63,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "enumerate files on Windows",
      "parse PE header",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "debug build",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [
        "45.33.23.183",
        "96.126.123.244"
      ],
      "domains": [
        "ftp.0adz.com",
        "ftp.search.instantfreesite.com",
        "gripto.leadhoster.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe"
      ],
      "registry_all_count": 5,
      "files_dropped": [
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "jusched.exe"
      ],
      "mutexes": [
        "RasPbFile",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\Program Files\\48db4bad\\jusched.exe",
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \""
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Rising MOVES",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}