# dataset/preprocessed/20cd2a370ca473a8b90f47eae467b940ad8c32b7694904ff6b8fa24bcc3bbef4.json
{
  "hash": "20cd2a370ca473a8b90f47eae467b940ad8c32b7694904ff6b8fa24bcc3bbef4",
  "report_text": "--- NETWORK ACTIVITY ---\n\nSERVICES OPENED (1):\n  Messenger\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  f5dd26e2883ae915635a2cb2cecf6bb8.virus.exe\n    C:\\Users\\<USER>\\Downloads\\f5dd26e2883ae915635a2cb2cecf6bb8.virus.exe\n    C:\\Users\\<USER>\\Downloads\\f5dd26e2883ae915635a2cb2cecf6bb8.virus.exe\n\nCOMMAND EXECUTIONS (1):\n  C:\\Users\\<USER>\\Downloads\\f5dd26e2883ae915635a2cb2cecf6bb8.virus.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (16):\n  C:\\Windows\\mssrv.exe  [sha256: 071b3e2b2e11ef5f...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\xxx [milf] lady .zip.exe  [sha256: 0540138931e9dc2b...]\n  C:\\Program Files\\dotnet\\shared\\russian gang bang xxx [bangbus] gorgeoushorny .avi.exe  [sha256: 49d220e4ee347648...]\n  C:\\Program Files\\DVD Maker\\Shared\\xxx hidden titts shoes .mpeg.exe  [sha256: b685447d5e61bceb...]\n  C:\\Program Files\\Windows Journal\\Templates\\bukkake uncut beautyfull .mpeg.exe  [sha256: 9cd1d5c935c8aa5f...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\hardcore sleeping titts stockings .zip.exe  [sha256: 7402b07c7acfe260...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\lingerie catfight sm .mpeg.exe  [sha256: 26d552b588c0d504...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\xxx full movie hotel .mpeg.exe  [sha256: c6592d09acdea5c7...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\japanese gang bang bukkake [bangbus] 40+ .avi.exe  [sha256: 1fa53ee11f511023...]\n  C:\\Program Files (x86)\\dotnet\\shared\\swedish handjob trambling catfight feet .avi.exe  [sha256: 94d109884dd893f9...]\n  C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\tyrkish cumshot blowjob catfight traffic .zip.exe  [sha256: 2bd67197c395f2ca...]\n  C:\\Program Files (x86)\\Google\\Update\\Download\\russian beastiality lesbian public titts blondie .mpeg.exe  [sha256: 8f3656abf9ded5ef...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\black animal trambling full movie fishy .zip.exe  [sha256: 46357801adba588c...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\brasilian action hardcore licking  (Karin).zip.exe  [sha256: bdb9e54d693629d9...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\indian fetish xxx public glans penetration  (Tatjana).avi.exe  [sha256: aa8b0e06f579493e...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\brasilian cumshot bukkake [free] glans girly .avi.exe  [sha256: f3443a3342c3d047...]\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880266189575195  size=432318\n",
  "ground_truth": {
    "hash": "20cd2a370ca473a8b90f47eae467b940ad8c32b7694904ff6b8fa24bcc3bbef4",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\lingerie catfight sm .mpeg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\xxx full movie hotel .mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\japanese gang bang bukkake [bangbus] 40+ .avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\tyrkish cumshot blowjob catfight traffic .zip.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\russian beastiality lesbian public titts blondie .mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\black animal trambling full movie fishy .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\brasilian action hardcore licking  (Karin).zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\brasilian cumshot bukkake [free] glans girly .avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\indian fetish xxx public glans penetration  (Tatjana).avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\swedish handjob trambling catfight feet .avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\xxx [milf] lady .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\xxx hidden titts shoes .mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\bukkake uncut beautyfull .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\hardcore sleeping titts stockings .zip.exe",
        "C:\\Program Files\\dotnet\\shared\\russian gang bang xxx [bangbus] gorgeoushorny .avi.exe",
        "C:\\Windows\\mssrv.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Jujubox"
  ]
}