# dataset/preprocessed/21b251de1ec869751bdbf7fa63769b25d97b1a73e5c30c9a88594f9cb5bc98f7.json
{
  "hash": "21b251de1ec869751bdbf7fa63769b25d97b1a73e5c30c9a88594f9cb5bc98f7",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  danish porn horse licking titts .rar.exe\n    C:\\Users\\<USER>\\Downloads\\danish porn horse licking titts .rar.exe\n    C:\\Users\\<USER>\\Downloads\\danish porn horse licking titts .rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\21b251de1ec869751bdbf7fa63769b25d97b1a73e5c30c9a88594f9cb5bc98f7(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\danish porn horse licking titts .rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american animal beast voyeur hairy  (Christine,Karin).avi.exe  [sha256: e13e5300d7007510...]\n  beast hidden sweet  (Sandy,Tatjana).avi.exe  [sha256: cf54dd0a9e4deeed...]\n  brasilian action hardcore catfight gorgeoushorny .avi.exe  [sha256: 35a470d2f5603e14...]\n  brasilian action lesbian several models upskirt .avi.exe  [sha256: fcb759e52402e591...]\n  bukkake big cock pregnant  (Sarah).rar.exe  [sha256: 7a87d80902fae97f...]\n  bukkake uncut penetration .avi.exe  [sha256: 111049ff5b3d5cb1...]\n  bukkake [milf] penetration .mpg.exe  [sha256: 161dfd8ef585c88b...]\n  danish beastiality blowjob [bangbus]  (Janette).mpeg.exe  [sha256: d0965119678dd47a...]\n  danish fetish horse licking .zip.exe  [sha256: 65d6536fb5e8185d...]\n  defrag.exe  [sha256: a670c974124e7e50...]\n  defrag64.exe  [sha256: cb1ed788d31b7d53...]\n  deinstall.exe  [sha256: 3c30f486b954c483...]\n  driver.exe  [sha256: 432bf5759ebf3fa0...]\n  drwatson32.exe  [sha256: acacf43237d45cfd...]\n  fucking hidden titts boots  (Curtney).zip.exe  [sha256: c66de768eee17a9e...]\n  gay masturbation .avi.exe  [sha256: c04a5017127c2db0...]\n  hardcore big shoes  (Sonja,Sylvia).mpg.exe  [sha256: 5398d0a92f5ca45b...]\n  IE9setup.exe  [sha256: 9d9c247b5daf345d...]\n  indian cum xxx hidden hole  (Jenna,Janette).avi.exe  [sha256: e881826b5ae6b409...]\n  italian fetish hardcore [free] redhair .rar.exe  [sha256: 984498bb072febd1...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\gay masturbation .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\tyrkish animal sperm hidden feet .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\italian porn sperm uncut beautyfull  (Kathrin,Samantha).zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\italian kicking bukkake licking cock mistress .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\russian cum bukkake lesbian 40+ .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\swedish fetish trambling uncut .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\tyrkish animal beast girls feet circumcision  (Liz).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\russian kicking lingerie voyeur gorgeoushorny .rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\fucking hidden titts boots  (Curtney).zip.exe\n  C:\\Program Files\\Movie Maker\\Shared\\japanese porn xxx [free] glans  (Ashley,Tatjana).avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\russian kicking blowjob [bangbus] glans .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\japanese fetish lesbian hidden traffic  (Sandy,Karin).avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\hardcore big shoes  (Sonja,Sylvia).mpg.exe\n  C:\\WINDOWS\\Temp\\beast hidden sweet  (Sandy,Tatjana).avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\russian nude lingerie several models sweet  (Christine,Jade).mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\danish fetish horse licking .zip.exe\n  C:\\WINDOWS\\assembly\\temp\\russian beastiality lingerie [free] .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\bukkake uncut penetration .avi.exe\n  C:\\WINDOWS\\ime\\shared\\bukkake big cock pregnant  (Sarah).rar.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880060195922852  size=1907018\n",
  "ground_truth": {
    "hash": "21b251de1ec869751bdbf7fa63769b25d97b1a73e5c30c9a88594f9cb5bc98f7",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\xxx uncut cock mature  (Sarah).mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\hardcore [free] balls .mpg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\lingerie masturbation .mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\horse licking ash .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\indian handjob blowjob big stockings .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\swedish cumshot trambling catfight  (Tatjana).avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\tyrkish nude bukkake catfight hole femdom .rar.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\tyrkish porn trambling [free] stockings .rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\horse hardcore full movie hole granny  (Janette).rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\japanese handjob gay voyeur 50+ .mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\sperm hot (!) .mpeg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american animal beast voyeur hairy  (Christine,Karin).avi.exe",
        "beast hidden sweet  (Sandy,Tatjana).avi.exe",
        "brasilian action hardcore catfight gorgeoushorny .avi.exe",
        "brasilian action lesbian several models upskirt .avi.exe",
        "bukkake [milf] penetration .mpg.exe",
        "bukkake big cock pregnant  (Sarah).rar.exe",
        "bukkake uncut penetration .avi.exe",
        "danish beastiality blowjob [bangbus]  (Janette).mpeg.exe",
        "danish fetish horse licking .zip.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking hidden titts boots  (Curtney).zip.exe",
        "gay masturbation .avi.exe",
        "hardcore big shoes  (Sonja,Sylvia).mpg.exe",
        "indian cum xxx hidden hole  (Jenna,Janette).avi.exe",
        "italian fetish hardcore [free] redhair .rar.exe",
        "italian kicking bukkake licking cock mistress .mpg.exe",
        "italian porn sperm uncut beautyfull  (Kathrin,Samantha).zip.exe",
        "japanese fetish lesbian hidden traffic  (Sandy,Karin).avi.exe",
        "japanese porn xxx [free] glans  (Ashley,Tatjana).avi.exe",
        "mssrv.exe",
        "russian beastiality lingerie [free] .avi.exe",
        "russian cum bukkake lesbian 40+ .mpeg.exe",
        "russian kicking blowjob [bangbus] glans .mpg.exe",
        "russian kicking lingerie voyeur gorgeoushorny .rar.exe",
        "russian nude lingerie several models sweet  (Christine,Jade).mpg.exe",
        "screensaver.scr",
        "soccer.exe",
        "start.exe",
        "swedish fetish trambling uncut .mpeg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}