# dataset/preprocessed/2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb400975f.json
{
  "hash": "2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb400975f",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  2336c490685cc56e_program files .exe\n    C:\\ProgramData\\roabdd.exe\n\nCOMMAND EXECUTIONS (2):\n  (open) C:\\2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb40(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\2336c490685cc56e_program files\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - contain obfuscated stackstrings\n  - get startup folder\n  - set file attributes\n  - get common file path\n  - accept command line arguments\n\nCAPA RULE MATCHES (MBC):\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - get startup folder\n  - Process::Create Process [C0017]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - contain loop\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Delete File [C0047]\n  - File System::Create Directory [C0046]\n  - Process::Allocate Thread Local Storage [C0040]\n  - get thread local storage value\n  - Process::Set Thread Local Storage Value [C0041]\n  - File System::Read File [C0051]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - resolve function by parsing PE exports\n  - contain a resource (.rsrc) section\n  - Discovery::Code Discovery::Enumerate PE Sections [B0046.001]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb400975f .exe  [sha256: 8c6a3644808c93e8...]\n  AUTOEXEC.BAT .exe  [sha256: a86d523c98168cc0...]\n  boot.ini .exe  [sha256: 569b7a538be8ade1...]\n  CONFIG.SYS .exe  [sha256: a8839e78ea244b3d...]\n  cuckoo .exe  [sha256: c00f15decb16d732...]\n  Documents and Settings .exe  [sha256: d2f47df9a1b5de84...]\n  eula.1028.txt .exe  [sha256: 85f5a95ca5b7e798...]\n  eula.1031.txt .exe  [sha256: 23d37c707b068c42...]\n  eula.1033.txt .exe  [sha256: b2ecae7864772410...]\n  eula.1036.txt .exe  [sha256: 90fc4e42b2fb83e4...]\n  eula.1040.txt .exe  [sha256: 8186ebfa0babb439...]\n  eula.1041.txt .exe  [sha256: 84f222e0f71cbe8d...]\n  eula.1042.txt .exe  [sha256: 8940714a1b771acb...]\n  eula.2052.txt .exe  [sha256: 9af71f2199484448...]\n  eula.3082.txt .exe  [sha256: c15d95c183b415d6...]\n  globdata.ini .exe  [sha256: 79b97f544730c45e...]\n  install.exe .exe  [sha256: 61c5e462fd66de0e...]\n  install.ini .exe  [sha256: b07cfbb438959c42...]\n  install.res.1028.dll .exe  [sha256: a8202f6a978c6782...]\n  install.res.1031.dll .exe  [sha256: 3ab9312ea1d3e2b2...]\n\nFILES WRITTEN (30):\n  C:\\0 .exe\n  C:\\AUTOEXEC.BAT .exe\n  C:\\CONFIG.SYS .exe\n  C:\\DiskD .exe\n  C:\\DiskX .exe\n  C:\\DllLoader.exe .exe\n  C:\\Documents and Settings .exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Saaaalamm\\Mira.h\n  C:\\Documents and Settings\\All Users\\Application Data\\qbqvqe.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\vypfkm.exe\n  C:\\ProgramData\\Saaaalamm\\Mira.h\n  C:\\ProgramData\\roabdd.exe\n  C:\\boot.ini .exe\n  C:\\bootfont.bin .exe\n  C:\\eula.2052.txt .exe\n  C:\\globdata.ini .exe\n  C:\\hiberfil.sys .exe\n  C:\\hosts .exe\n  c:\\$Recycle.Bin .exe\n  c:\\2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb400975f .exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft\u00ae Windows\u00ae Operating System = C:\\ProgramData\\roabdd.exe\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System = C:\\Documents and Settings\\All Users\\Application Data\\qbqvqe.exe\u0000\n\nMUTEXES CREATED (7):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\All Users\\Application Data\\vypfkm.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA\n\nCOMPILER AND PACKER:\n  [Compiler] MinGW\n  [Linker] GNU linker ld (GNU Binutils)\n  Packers: {'PEiD': 'MingWin32 v?.? (h)'}\n\nPE SECTIONS (entropy):\n  .text        entropy=5.78\n  .rsrc        entropy=6.88\n  [overlay]    entropy=6.950550079345703  size=388284\n",
  "ground_truth": {
    "hash": "2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb400975f",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8714,
    "malicious_count": 61,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "get startup folder"
      }
    ],
    "capa_sig_names": [
      "contain obfuscated stackstrings",
      "get startup folder",
      "create process on Windows",
      "set file attributes",
      "delay execution",
      "create or open file",
      "contain loop",
      "get common file path",
      "delete file",
      "create directory",
      "allocate thread local storage",
      "get thread local storage value",
      "set thread local storage value",
      "read file on Windows",
      "accept command line arguments",
      "resolve function by parsing PE exports",
      "contain a resource (.rsrc) section",
      "create or open file",
      "get common file path",
      "read file on Windows",
      "enumerate PE sections"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft\u00ae Windows\u00ae Operating System",
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "2336c490685cc56e38ceaccf11480c616ece18fc27a1abc232d2025bb400975f .exe",
        "AUTOEXEC.BAT .exe",
        "C:\\ProgramData\\Saaaalamm\\Mira.h",
        "C:\\ProgramData\\roabdd.exe",
        "CONFIG.SYS .exe",
        "Documents and Settings .exe",
        "IO.SYS .exe",
        "MSDOS.SYS .exe",
        "Mira.h",
        "NTDETECT.COM .exe",
        "Program Files .exe",
        "Python27 .exe",
        "RECYCLER .exe",
        "System Volume Information .exe",
        "VC_RED.MSI .exe",
        "VC_RED.cab .exe",
        "WINDOWS .exe",
        "boot.ini .exe",
        "c:\\$Recycle.Bin .exe",
        "c:\\Documents and Settings .exe",
        "c:\\PerfLogs .exe",
        "c:\\Program Files (x86) .exe",
        "c:\\Program Files .exe",
        "c:\\ProgramData .exe",
        "c:\\dev_man.exe .exe",
        "c:\\l.bat .exe",
        "c:\\pagefile.sys .exe",
        "c:\\r.py .exe",
        "cuckoo .exe",
        "eula.1028.txt .exe",
        "eula.1031.txt .exe",
        "eula.1033.txt .exe",
        "eula.1036.txt .exe",
        "eula.1040.txt .exe",
        "eula.1041.txt .exe",
        "eula.1042.txt .exe",
        "eula.2052.txt .exe",
        "eula.3082.txt .exe",
        "globdata.ini .exe",
        "install.exe .exe",
        "install.ini .exe",
        "install.res.1028.dll .exe",
        "install.res.1031.dll .exe",
        "install.res.1033.dll .exe",
        "install.res.1036.dll .exe",
        "install.res.1040.dll .exe",
        "install.res.1041.dll .exe",
        "install.res.1042.dll .exe",
        "install.res.2052.dll .exe",
        "install.res.3082.dll .exe",
        "ntldr .exe",
        "pagefile.sys .exe",
        "vcredist.bmp .exe",
        "vypfkm.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\Documents and Settings\\All Users\\Application Data\\vypfkm.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}