# dataset/preprocessed/23bb34c10354cdc1565f5a7eb48fdc633ade7dc8ba52e187440aa4005fadc481.json
{
  "hash": "23bb34c10354cdc1565f5a7eb48fdc633ade7dc8ba52e187440aa4005fadc481",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  deinstall.exe\n    C:\\Users\\<USER>\\Downloads\\deinstall.exe\n    C:\\Users\\<USER>\\Downloads\\deinstall.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\23bb34c10354cdc1565f5a7eb48fdc633ade7dc8ba52e187440aa4005fadc481(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\deinstall.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  action cum catfight .zip.exe  [sha256: 0534482dff1cdf11...]\n  asian animal full movie latex .rar.exe  [sha256: cdd851e63868f958...]\n  asian blowjob full movie sm .avi.exe  [sha256: e21db0d329f0b150...]\n  asian cumshot lesbian cock pregnant .mpg.exe  [sha256: c1feffc865149291...]\n  beast cum lesbian 50+ .mpg.exe  [sha256: c6079fac419ccea0...]\n  british gang bang xxx [free] sm .mpeg.exe  [sha256: f248d6afdb2ffb41...]\n  bukkake [milf] .mpeg.exe  [sha256: ce5ec2a4f84121d6...]\n  canadian handjob full movie shoes .mpg.exe  [sha256: df471e3be8da7fa7...]\n  chinese fucking blowjob licking feet black hairunshaved .rar.exe  [sha256: 66760fcaa48a23db...]\n  chinese sperm bukkake [bangbus] cock blondie .mpeg.exe  [sha256: 6ab577c1a8f128fc...]\n  cumshot catfight .zip.exe  [sha256: 3bc341567b8faf63...]\n  danish handjob trambling catfight blondie .mpg.exe  [sha256: 1c16fdf3f60d7f60...]\n  danish hardcore [milf] .avi.exe  [sha256: a6e5dc7679062872...]\n  danish nude xxx lesbian sm  (Curtney,Janette).rar.exe  [sha256: 32d4332f1a0dbbbb...]\n  defrag.exe  [sha256: 8c0e7f818f490b13...]\n  deinstall.exe  [sha256: 181b23d5c273248a...]\n  driver.exe  [sha256: 21813de425fb6bfa...]\n  drwatson32.exe  [sha256: fe22c22b0ea24821...]\n  fetish catfight hole penetration .rar.exe  [sha256: 6dbf1eb8057ea7c5...]\n  fetish several models sweet  (Tatjana,Jade).avi.exe  [sha256: 2cd2876c1fa7c5fe...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\japanese horse lingerie uncut .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\british gang bang xxx [free] sm .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\asian blowjob full movie sm .avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\danish hardcore [milf] .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\indian nude licking traffic .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\german animal uncut nipples  (Sonja).mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\norwegian gay masturbation hole  (Ashley).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\danish handjob trambling catfight blondie .mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\indian beastiality lingerie uncut titts .rar.exe\n  C:\\Program Files\\Movie Maker\\Shared\\chinese fucking blowjob licking feet black hairunshaved .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\action cum catfight .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\trambling uncut mistress  (Christine).mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\fetish several models sweet  (Tatjana,Jade).avi.exe\n  C:\\WINDOWS\\Temp\\asian animal full movie latex .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\hardcore lingerie hidden pregnant .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\fetish catfight hole penetration .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\malaysia blowjob kicking public wifey .zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\chinese sperm bukkake [bangbus] cock blondie .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\indian handjob voyeur nipples .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.879809856414795  size=307404\n",
  "ground_truth": {
    "hash": "23bb34c10354cdc1565f5a7eb48fdc633ade7dc8ba52e187440aa4005fadc481",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "action cum catfight .zip.exe",
        "asian animal full movie latex .rar.exe",
        "asian blowjob full movie sm .avi.exe",
        "asian cumshot lesbian cock pregnant .mpg.exe",
        "beast cum lesbian 50+ .mpg.exe",
        "british gang bang xxx [free] sm .mpeg.exe",
        "bukkake [milf] .mpeg.exe",
        "canadian handjob full movie shoes .mpg.exe",
        "chinese fucking blowjob licking feet black hairunshaved .rar.exe",
        "chinese sperm bukkake [bangbus] cock blondie .mpeg.exe",
        "cumshot catfight .zip.exe",
        "danish handjob trambling catfight blondie .mpg.exe",
        "danish hardcore [milf] .avi.exe",
        "danish nude xxx lesbian sm  (Curtney,Janette).rar.exe",
        "defrag.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "fetish catfight hole penetration .rar.exe",
        "fetish several models sweet  (Tatjana,Jade).avi.exe",
        "french trambling voyeur feet .rar.exe",
        "game3d.exe",
        "german animal uncut nipples  (Sonja).mpg.exe",
        "german horse action lesbian 50+ .rar.exe",
        "german nude full movie glans .mpg.exe",
        "hardcore lingerie hidden pregnant .mpg.exe",
        "indian beastiality lingerie uncut titts .rar.exe",
        "indian handjob voyeur nipples .zip.exe",
        "indian nude licking traffic .mpeg.exe",
        "install32.exe",
        "japanese horse lingerie uncut .mpg.exe",
        "malaysia blowjob kicking public wifey .zip.exe",
        "mssrv.exe",
        "norwegian gay masturbation hole  (Ashley).zip.exe",
        "notes.txt.exe",
        "pgp9.exe",
        "safetyserver.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}