# dataset/preprocessed/247a377514ac91db421f5a93a6f7d42d06636232a458edeb89bd0727cd465134.json
{
  "hash": "247a377514ac91db421f5a93a6f7d42d06636232a458edeb89bd0727cd465134",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ****.exe\n      ****.exe\n      ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n    ****.exe\n\nCOMMAND EXECUTIONS (4):\n  ((null)) C:\\247a377514ac91db421f5a93a6f7d42d06636232a458edeb89bd0727cd465134(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\247a377514ac91db_black fetish hardcore uncut feet lady  (karin).mpg.exe\n  \"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\" \n  \"C:\\Windows\\system32\\rundll32.exe\" C:\\Windows\\system32\\shell32.dll,OpenAs_RunDLL C:\\Users\\Petra\\AppData\\Local\\Temp\\black-fetish-hardcor.pe32\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american gang bang horse catfight penetration .rar.exe  [sha256: 280a21021cf87fc5...]\n  american porn xxx big glans  (Britney,Sylvia).mpeg.exe  [sha256: 3a9de7665752f9cf...]\n  blowjob [free] .zip.exe  [sha256: f75e9c65980e1298...]\n  brasilian animal sperm full movie mature .avi.exe  [sha256: 8655664ecbef3309...]\n  brasilian kicking xxx big balls .zip.exe  [sha256: 8f7d4ff05a7a8ec9...]\n  bukkake public penetration .mpeg.exe  [sha256: 459968c546df738f...]\n  danish fetish beast public titts granny  (Janette).mpeg.exe  [sha256: c70628005f6b3299...]\n  defrag.exe  [sha256: 1102adb415edf443...]\n  defrag64.exe  [sha256: 8794ab16441c8fe8...]\n  deinstall.exe  [sha256: b82daaa0c2efc9bf...]\n  driver.exe  [sha256: 24baa11523fe16a8...]\n  fucking catfight 40+ .rar.exe  [sha256: e6380cd98de75c5f...]\n  fucking sleeping  (Karin).zip.exe  [sha256: c390933173e09e1d...]\n  gay catfight glans .mpeg.exe  [sha256: 2e76986a4b9e0828...]\n  gay lesbian wifey .avi.exe  [sha256: 5eb8b3e86b1180c9...]\n  gay licking swallow  (Britney,Tatjana).mpeg.exe  [sha256: 3d1539a0746bf050...]\n  gay [bangbus] titts bedroom .mpeg.exe  [sha256: b6c43793c2ee36f9...]\n  horse catfight mature .mpg.exe  [sha256: c95588294feba865...]\n  horse hot (!) .mpeg.exe  [sha256: 4367011374d51582...]\n  IE9setup.exe  [sha256: 6f27c5f672fed99b...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\russian animal fucking [free] .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\blowjob [free] .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\indian beastiality lesbian lesbian .mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\lesbian several models .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\xxx several models blondie .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\italian animal bukkake hot (!) glans .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\xxx uncut glans .mpg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\gay [bangbus] titts bedroom .mpeg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\russian fetish gay voyeur redhair .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\horse hot (!) .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\sperm sleeping glans .mpeg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\gay catfight glans .mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\brasilian animal sperm full movie mature .avi.exe\n  C:\\WINDOWS\\Temp\\sperm hot (!) young .mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\tyrkish horse trambling catfight circumcision .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\swedish cumshot xxx hot (!)  (Karin).zip.exe\n  C:\\WINDOWS\\assembly\\temp\\american gang bang horse catfight penetration .rar.exe\n  C:\\WINDOWS\\assembly\\tmp\\indian fetish xxx sleeping .zip.exe\n  C:\\WINDOWS\\ime\\shared\\fucking sleeping  (Karin).zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\WINDOWS\\mssrv.exe\u0000\n\nOTHER REGISTRY KEYS SET (2):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n  \\REGISTRY\\MACHINE\\HexStream = 53004F004600540057004100520045005C004D006900630072006F0073006F00660074005C00570069006E0064\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n  mutex666\n\nPROCESSES CREATED (2):\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe\n  rundll32.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.64\n  [overlay]    entropy=7.879457473754883  size=204993\n",
  "ground_truth": {
    "hash": "247a377514ac91db421f5a93a6f7d42d06636232a458edeb89bd0727cd465134",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 4,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\asian gay animal catfight  (Christine,Sandy).mpg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\black cum big hole .mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\french nude xxx [milf] sm .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\japanese trambling [bangbus] girly .rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\russian handjob gay voyeur black hairunshaved .mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\tyrkish lingerie animal girls .zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\xxx xxx uncut .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\beast handjob voyeur ejaculation  (Sylvia,Curtney).avi.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "\\Program Files\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\hardcore animal public  (Sylvia).zip.exe",
        "\\Program Files\\Common Files\\microsoft shared\\british handjob voyeur ash .rar.exe",
        "\\Windows\\mssrv.exe",
        "american gang bang horse catfight penetration .rar.exe",
        "american porn xxx big glans  (Britney,Sylvia).mpeg.exe",
        "blowjob [free] .zip.exe",
        "brasilian animal sperm full movie mature .avi.exe",
        "brasilian kicking xxx big balls .zip.exe",
        "bukkake public penetration .mpeg.exe",
        "danish fetish beast public titts granny  (Janette).mpeg.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "fucking catfight 40+ .rar.exe",
        "fucking sleeping  (Karin).zip.exe",
        "gay [bangbus] titts bedroom .mpeg.exe",
        "gay catfight glans .mpeg.exe",
        "gay lesbian wifey .avi.exe",
        "gay licking swallow  (Britney,Tatjana).mpeg.exe",
        "horse catfight mature .mpg.exe",
        "horse hot (!) .mpeg.exe",
        "indian beastiality lesbian lesbian .mpg.exe",
        "indian fetish xxx sleeping .zip.exe",
        "italian animal bukkake hot (!) glans .zip.exe",
        "italian horse lingerie full movie circumcision .avi.exe",
        "lesbian several models .mpeg.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian animal fucking [free] .mpeg.exe",
        "russian fetish bukkake several models  (Sylvia).mpeg.exe",
        "russian fetish gay voyeur redhair .mpg.exe",
        "screensaver.scr",
        "soccer.exe",
        "sperm hot (!) young .mpg.exe",
        "sperm sleeping glans .mpeg.exe",
        "swedish cumshot xxx hot (!)  (Karin).zip.exe",
        "tyrkish horse trambling catfight circumcision .mpg.exe",
        "undelete.exe",
        "xxx several models blondie .mpg.exe",
        "xxx uncut glans .mpg.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex",
        "mutex666"
      ],
      "processes": [
        "C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\EB93A6\\996E.exe",
        "rundll32.exe"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "SNDBOX",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}