# dataset/preprocessed/250b39dbdd193c0f5464f7d1b9922b46b32754d957a397874c795e75ee41bb82.json
{
  "hash": "250b39dbdd193c0f5464f7d1b9922b46b32754d957a397874c795e75ee41bb82",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  C:\\Program Files\\DVD Maker\\Shared\\indian cum [free] vagina upskirt  (Liz,Gina).rar.exe\n    C:\\Program Files\\DVD Maker\\Shared\\indian cum [free] vagina upskirt  (Liz,Gina).rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\250b39dbdd193c0f5464f7d1b9922b46b32754d957a397874c795e75ee41bb82(null) [(null)]\n  C:\\Program Files\\DVD Maker\\Shared\\indian cum [free] vagina upskirt  (Liz,Gina).rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american gang bang xxx catfight .mpeg.exe  [sha256: d7f4c64c503e5c7a...]\n  american nude sperm public  (Jade).rar.exe  [sha256: efadf4950e50a0fb...]\n  american porn lingerie [bangbus] lady .zip.exe  [sha256: a09aa5aa0d072d8f...]\n  beast catfight pregnant .rar.exe  [sha256: 7838f05d2ff3d7ae...]\n  beast hot (!) glans granny  (Jade).zip.exe  [sha256: a6108a37b47653c8...]\n  beast voyeur .avi.exe  [sha256: e466ccfddba2b396...]\n  black horse sperm [free] .mpg.exe  [sha256: bebdbff6e2ebd0d1...]\n  black kicking horse hidden .zip.exe  [sha256: d00eaec23be09a2c...]\n  brasilian beastiality lingerie hot (!) feet .mpeg.exe  [sha256: 7a81035eff270ddd...]\n  brasilian horse blowjob full movie feet .avi.exe  [sha256: 8e1366dba9a92430...]\n  british trambling full movie glans shoes .rar.exe  [sha256: bbad7bb325aa42c1...]\n  defrag.exe  [sha256: 3fc02afc4a0ef87a...]\n  defrag64.exe  [sha256: 2aa0b964aa1d4fe8...]\n  deinstall.exe  [sha256: 46febfce4fa10f01...]\n  driver.exe  [sha256: 891557ac55c1735b...]\n  drwatson32.exe  [sha256: a2de328cb3f3476f...]\n  fucking big shoes .zip.exe  [sha256: 675fe18e4be471b8...]\n  fucking full movie feet .zip.exe  [sha256: f2f193f75e0b014f...]\n  fucking [bangbus] black hairunshaved .rar.exe  [sha256: 6c12f729ac650a48...]\n  horse big .zip.exe  [sha256: 263fdfd3fc5b3615...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\black kicking horse hidden .zip.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\italian animal lingerie full movie hairy .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\beast hot (!) glans granny  (Jade).zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\beast voyeur .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\fucking full movie feet .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\horse big .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\tyrkish cumshot beast uncut titts beautyfull .mpg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\lesbian [bangbus] hole sm  (Tatjana).mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\brasilian beastiality lingerie hot (!) feet .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\tyrkish gang bang lesbian [free] feet .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\british trambling full movie glans shoes .rar.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\beast catfight pregnant .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\lesbian uncut .mpeg.exe\n  C:\\WINDOWS\\Temp\\italian beastiality gay full movie gorgeoushorny .avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\trambling full movie cock 40+ .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\american porn lingerie [bangbus] lady .zip.exe\n  C:\\WINDOWS\\assembly\\temp\\tyrkish horse sperm catfight boots  (Gina,Jade).mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\fucking big shoes .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\fucking [bangbus] black hairunshaved .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.879866123199463  size=492344\n",
  "ground_truth": {
    "hash": "250b39dbdd193c0f5464f7d1b9922b46b32754d957a397874c795e75ee41bb82",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\lesbian masturbation  (Karin).rar.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\lesbian licking swallow .zip.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\animal hidden black hairunshaved  (Gina,Samantha).avi.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\german beastiality trambling [bangbus] ejaculation  (Britney).mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\handjob licking nipples wifey .avi.exe",
        "C:\\Program Files\\dotnet\\shared\\italian cumshot licking .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american gang bang xxx catfight .mpeg.exe",
        "american nude sperm public  (Jade).rar.exe",
        "american porn lingerie [bangbus] lady .zip.exe",
        "beast catfight pregnant .rar.exe",
        "beast hot (!) glans granny  (Jade).zip.exe",
        "beast voyeur .avi.exe",
        "black horse sperm [free] .mpg.exe",
        "black kicking horse hidden .zip.exe",
        "brasilian beastiality lingerie hot (!) feet .mpeg.exe",
        "brasilian horse blowjob full movie feet .avi.exe",
        "british trambling full movie glans shoes .rar.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking [bangbus] black hairunshaved .rar.exe",
        "fucking big shoes .zip.exe",
        "fucking full movie feet .zip.exe",
        "horse [free] swallow .rar.exe",
        "horse big .zip.exe",
        "horse public latex .mpeg.exe",
        "italian animal lingerie full movie hairy .zip.exe",
        "italian beastiality gay full movie gorgeoushorny .avi.exe",
        "lesbian [bangbus] hole sm  (Tatjana).mpg.exe",
        "lesbian uncut .mpeg.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian porn horse hot (!) glans latex .mpeg.exe",
        "screensaver.scr",
        "setup.exe",
        "soccer.exe",
        "swedish cum lingerie public .avi.exe",
        "trambling full movie cock 40+ .mpeg.exe",
        "trambling masturbation bondage .mpeg.exe",
        "tyrkish cumshot beast uncut titts beautyfull .mpg.exe",
        "tyrkish gang bang lesbian [free] feet .mpeg.exe",
        "tyrkish horse sperm catfight boots  (Gina,Jade).mpeg.exe",
        "undelete.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}