# dataset/preprocessed/25226d7e7728b2fabce7d6f80aee2ca2dbbcaf9d773f083cd59e462180e01bca.json
{
  "hash": "25226d7e7728b2fabce7d6f80aee2ca2dbbcaf9d773f083cd59e462180e01bca",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (1):\n  GET http://www.universal101.com/upd.sc  UA: IExplorer 34\n\nDNS LOOKUPS (1):\n  www.universal101.com \u2192 76.223.54.146, 13.248.169.48\n\nIP TRAFFIC (2):\n  114.114.114.114:53 (UDP)\n  13.248.169.48:80 (TCP)\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  25226d7e7728b2fa_divx32.dll.exe\n    C:\\Windows\\system32\\lssmon.exe\n      C:\\Windows\\system32\\srtsrv32.exe\n      C:\\Windows\\system32\\srtsrv32.exe\n      C:\\Windows\\system32\\srtsrv32.exe\n        C:\\Windows\\system32\\LSASSMGR.EXE\n    C:\\Windows\\system32\\srtsrv32.exe\n      C:\\Windows\\system32\\LSASSMGR.EXE\n        C:\\Windows\\system32\\LSASSMGR.EXE\n\nCOMMAND EXECUTIONS (3):\n  C:\\Windows\\system32\\srtsrv32.exe\n  C:\\Windows\\system32\\LSASSMGR.EXE\n  C:\\Windows\\system32\\lssmon.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (27):\n  - encode data using XOR\n  - contain obfuscated stackstrings\n  - hide graphical window\n  - access PEB ldr_data\n  - accept command line arguments\n  - query or enumerate registry value\n  - link function at runtime on Windows\n  - get geographical location\n  - find graphical window\n  - encrypt data using RC4 PRGA\n  - check if file exists\n  - get disk size\n  - get file version info\n  - query or enumerate registry key\n  - delete registry key\n  - delete registry value\n  - read clipboard data\n  - capture screenshot\n  - link many functions at runtime\n  - open clipboard\n  - log keystrokes via polling\n  - log keystrokes\n  - get keyboard layout\n  - get number of processors\n  - enumerate gui resources\n  - get hostname\n  - persist via Run registry key\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Process::Create Thread [C0038]\n  - File System::Create File [C0016]\n  - File System::Get File Attributes [C0049]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - hide graphical window\n  - set global application hook\n  - inspect load icon resource\n  - access PEB ldr_data\n  - Process::Create Process [C0017]\n  - contain loop\n  - Execution::Command and Scripting Interpreter [E1059]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - File System::Writes File [C0052]\n  - link function at runtime on Windows\n  - get geographical location\n\n--- ARTIFACTS ---\n\nFILES DROPPED (4):\n  C:\\Windows\\divx32.dll  [sha256: 20ce49628da4cb87...]\n  c:\\program files\\internet explorer\\iexplor.exe  [sha256: 8d69588880eefdb9...]\n  c:\\windows\\system32\\lssmon.exe  [sha256: 5d967727decd7c34...]\n  c:\\windows\\divx32.dll  [sha256: 5294b6aa8ca868bd...]\n\nFILES WRITTEN (6):\n  C:\\Program Files\\Internet Explorer\\iexplor.exe\n  C:\\WINDOWS\\system32\\LSASSMGR.EXE\n  C:\\WINDOWS\\system32\\spool.exe\n  C:\\WINDOWS\\system32\\srtsrv32.exe\n  C:\\Windows\\system32\\lssmon.exe\n  C:\\Windows\\system32\\srtsrv32.exe\n\nPERSISTENCE REGISTRY KEYS (12):\n  HKLM/SOFTWARE/Microsoft/Windows NT/CurrentVersion/Image File Execution Options//firefox.exe/Debugger = C:\\Program Files\\Mozilla Firefox\\firefoxe.exe\n  HKLM/SOFTWARE/Microsoft/Windows NT/CurrentVersion/Image File Execution Options//iexplore.exe/Debugger = C:\\Program Files\\Internet Explorer\\iexplor.exe\n  HKLM/SOFTWARE/Microsoft/Windows NT/CurrentVersion/Image File Execution Options//spoolsv.exe/Debugger = C:\\Windows\\system32\\spool.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\firefox.exe\\Debugger = C:\\Program Files (x86)\\Mozilla Firefox\\firefoxe.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\iexplore.exe\\Debugger = C:\\Program Files (x86)\\Internet Explorer\\iexplor.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\spoolsv.exe\\Debugger = C:\\Windows\\system32\\spool.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Layersecurity Servicemonitor = C:\\Windows\\system32\\lssmon.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\firefox.exe\\Debugger = C:\\Program Files\\Mozilla Firefox\\firefoxe.exe\u0000\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\iexplore.exe\\Debugger = C:\\Program Files\\Internet Explorer\\iexplor.exe\u0000\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\spoolsv.exe\\Debugger = C:\\WINDOWS\\system32\\spool.exe\u0000\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Layersecurity Servicemonitor = C:\\WINDOWS\\system32\\LSSMON.EXE\u0000\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Layersecurity Servicemonitor = C:\\WINDOWS\\system32\\lssmon.exe\u0000\n\nOTHER REGISTRY KEYS SET (10):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\WINDOWS\\\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\WINDOWS\\\n  HKU/S-1-5-21-470376811-3006406624-3672060426-1000/Software/Microsoft/Windows/CurrentVersion/Internet Settings/ZoneMap//U\n  HKLM/SOFTWARE/Microsoft/Windows/CurrentVersion/Run/Layersecurity Servicemonitor\n  HKLM/Software/Microsoft/Tracing/lssmon_RASMANCS/ConsoleTracingMask\n  HKLM/Software/Microsoft/Tracing/lssmon_RASMANCS/FileTracingMask\n  HKU/S-1-5-21-470376811-3006406624-3672060426-1000/Software/Microsoft/windows/CurrentVersion/Internet Settings/ProxyEnabl\n  HKU/S-1-5-21-470376811-3006406624-3672060426-1000/Software/Microsoft/Windows/CurrentVersion/Internet Settings/ZoneMap//A\n  HKLM/Software/Microsoft/Tracing/lssmon_RASAPI32/ConsoleTracingMask\n  HKLM/Software/Microsoft/Tracing/lssmon_RASAPI32/FileDirectory\n\nMUTEXES CREATED (25):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  IESQMMUTEX_0_208\n  LSM001\n  LSS001\n  Local\\!IETld!Mutex\n\nMUTEXES OPENED (2):\n  05387337392A20730D545F6985EB80691B51F7690E5FFF65803ABB6E015474081D15482073A5C76715\n  LSSUPD\n\nPROCESSES CREATED (4):\n  C:\\25226d7e7728b2fabce7d6f80aee2ca2dbbcaf9d773f083cd59e462180e01bca.exe\n  C:\\Windows\\system32\\LSASSMGR.EXE \n  C:\\Windows\\system32\\lssmon.exe \n  C:\\Windows\\system32\\srtsrv32.exe \n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  kernel32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [sfx] temporary EXE SFX\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  [Overlay] EXE file\n\nPE SECTIONS (entropy):\n  .text        entropy=6.59\n  .itext       entropy=6.36\n  .data        entropy=4.32\n  .bss         entropy=0.0\n  .idata       entropy=5.04\n  .tls         entropy=0.0\n  .rdata       entropy=0.21\n  .rsrc        entropy=6.78\n  [overlay]    entropy=6.703250408172607  size=2224128\n",
  "ground_truth": {
    "hash": "25226d7e7728b2fabce7d6f80aee2ca2dbbcaf9d773f083cd59e462180e01bca",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9167,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via polling"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "check if file exists"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry key"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1115",
        "signature_description": "read clipboard data"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "delay execution",
      "generate random numbers using the Delphi LCG",
      "create or open registry key",
      "create thread",
      "create or open file",
      "get file attributes",
      "PEB access",
      "allocate RWX memory",
      "encode data using XOR",
      "contain obfuscated stackstrings",
      "hide graphical window",
      "set global application hook",
      "inspect load icon resource",
      "access PEB ldr_data",
      "create process on Windows",
      "contain loop",
      "accept command line arguments",
      "query or enumerate registry value",
      "write file on Windows",
      "link function at runtime on Windows",
      "get geographical location",
      "set thread local storage value",
      "get thread local storage value",
      "create mutex",
      "find graphical window",
      "encrypt data using RC4 PRGA",
      "check if file exists",
      "get disk size",
      "get OS version",
      "get file version info",
      "extract resource via kernel32 functions",
      "query or enumerate registry key",
      "delete registry key",
      "delete registry value",
      "set registry value",
      "read clipboard data",
      "capture screenshot",
      "link many functions at runtime",
      "open clipboard",
      "write clipboard data",
      "replace clipboard data",
      "log keystrokes via polling",
      "log keystrokes",
      "get keyboard layout",
      "check for time delay via GetTickCount",
      "get number of processors",
      "enumerate gui resources",
      "set application hook",
      "get graphical window text",
      "get hostname",
      "check for software breakpoints",
      "create HTTP request",
      "connect to HTTP server",
      "copy file",
      "persist via Run registry key",
      "connect to URL",
      "read data from Internet",
      "receive data",
      "contain a thread local storage (.tls) section",
      "contain an embedded PE file",
      "contain a resource (.rsrc) section",
      "compiled with Borland Delphi",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [
        "114.114.114.114",
        "13.248.169.48"
      ],
      "domains": [
        "www.universal101.com"
      ],
      "urls": [
        "http://www.universal101.com/upd.sc"
      ],
      "registry_persistence": [
        "HKLM/SOFTWARE/Microsoft/Windows NT/CurrentVersion/Image File Execution Options//firefox.exe/Debugger",
        "HKLM/SOFTWARE/Microsoft/Windows NT/CurrentVersion/Image File Execution Options//iexplore.exe/Debugger",
        "HKLM/SOFTWARE/Microsoft/Windows NT/CurrentVersion/Image File Execution Options//spoolsv.exe/Debugger",
        "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\firefox.exe\\Debugger",
        "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\iexplore.exe\\Debugger",
        "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\spoolsv.exe\\Debugger",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Layersecurity Servicemonitor",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\firefox.exe\\Debugger",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\iexplore.exe\\Debugger",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Image File Execution Options\\spoolsv.exe\\Debugger",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Layersecurity Servicemonitor"
      ],
      "registry_all_count": 40,
      "files_dropped": [
        "C:\\Windows\\divx32.dll",
        "c:\\program files\\internet explorer\\iexplor.exe",
        "c:\\windows\\divx32.dll",
        "c:\\windows\\system32\\lssmon.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "IESQMMUTEX_0_208",
        "LSM001",
        "LSS001",
        "Local\\!IETld!Mutex",
        "Local\\WininetConnectionMutex",
        "Local\\WininetProxyRegistryMutex",
        "Local\\WininetStartupMutex",
        "Local\\ZoneAttributeCacheCounterMutex",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex",
        "Local\\_!MSFTHISTORY!_",
        "Local\\c:!users!admin!appdata!local!microsoft!windows!history!history.ie5!",
        "Local\\c:!users!admin!appdata!local!microsoft!windows!temporary internet files!content.ie5!",
        "Local\\c:!users!admin!appdata!roaming!microsoft!windows!cookies!",
        "RasPbFile",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "C:\\25226d7e7728b2fabce7d6f80aee2ca2dbbcaf9d773f083cd59e462180e01bca.exe",
        "C:\\Windows\\system32\\LSASSMGR.EXE ",
        "C:\\Windows\\system32\\lssmon.exe ",
        "C:\\Windows\\system32\\srtsrv32.exe "
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Sangfor ZSand",
    "Tencent HABO",
    "VirusTotal Jujubox"
  ]
}