# dataset/preprocessed/263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533.json
{
  "hash": "263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (6):\n  GET http://ctldl.windowsupdate.com/msdownload/update/v3/static/trustedr/en/disallowedcertstl.cab?53379364b6dd3cc4  UA: Microsoft-CryptoAPI/6.1\n  GET http://ocsp.pki.goog/gts1o1core/MFIwUDBOMEwwSjAJBgUrDgMCGgUABBRCRjDCJxnb3nDwj%2Fxz5aZfZjgXvAQUmNH4bhDrz5vsYJ8YkBug630J%2FSsCEQDM9j8%2B94pW6ggAAAAATZ8E  UA: Microsoft-CryptoAPI/6.1\n  GET http://ocsp.digicert.com/MFEwTzBNMEswSTAJBgUrDgMCGgUABBTBL0V27RVZ7LBduom%2FnYB45SPUEwQU5Z1ZMIJHWMys%2BghUNoZ7OrUETfACEAiIzVJfGSRETRSlgpHeuVI%3D  UA: Microsoft-CryptoAPI/6.1\n  CONNECT http://ctldl.windowsupdate.com/msdownload/update/v3/static/trustedr/en/disallowedcertstl.cab\n  GET https://www.google-analytics.com/collect?v=1&tid=UA-96780277-2&cid=75c4d8f9-d671-46d9-86b1-77e51760d218&t=event&aip=1&ec=Wallpaper_HD_MA&ea=Install&el=20180701&z=1597667476\n  GET https://wallhd.azurewebsites.net/Utils/CheckNewVersion?t=620&curVersion=1.4.0.12&target=newestSupported\n\nDNS LOOKUPS (13):\n  au-bg-shim.trafficmanager.net \u2192 67.26.131.254, 8.241.18.126\n  audownload.windowsupdate.nsatc.net \u2192 67.26.131.254, 8.241.18.126\n  auto.au.download.windowsupdate.com.c.footprint.net \u2192 67.26.131.254, 8.241.18.126\n  cs9.wac.phicdn.net \u2192 93.184.220.29\n  ctldl.windowsupdate.com \u2192 67.26.131.254, 8.241.18.126\n  ocsp.digicert.com \u2192 93.184.220.29\n  ocsp.pki.goog \u2192 172.217.18.163\n  pki-goog.l.google.com \u2192 172.217.18.163\n  wallhd.azurewebsites.net \u2192 52.178.105.179\n  waws-prod-am2-123.cloudapp.net \u2192 52.178.105.179\n  waws-prod-am2-123.vip.azurewebsites.windows.net \u2192 52.178.105.179\n  www-google-analytics.l.google.com \u2192 216.58.206.14\n  www.google-analytics.com \u2192 216.58.206.14\n\nIP TRAFFIC (5):\n  172.217.18.163:80 (TCP)\n  216.58.206.14:443 (TCP)\n  52.178.105.179:443 (TCP)\n  67.26.131.254:80 (TCP)\n  93.184.220.29:80 (TCP)\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533.exe\n    C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHd.exe /install 263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533.exe\n      C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdHelper64.exe /install\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (24):\n  - set file attributes\n  - get file system object information\n  - encode data using XOR\n  - query environment variable\n  - find graphical window\n  - enumerate files on Windows\n  - link function at runtime on Windows\n  - query or enumerate registry value\n  - query or enumerate registry key\n  - create shortcut via IShellLink\n  - delete registry value\n  - delete registry key\n  - get file size\n  - shutdown system\n  - get common file path\n  - accept command line arguments\n  - modify access privileges\n  - check OS version\n  - check if file exists\n  - get disk size\n  - capture webcam image\n  - open clipboard\n  - enumerate files recursively\n  - reference anti-VM strings targeting Xen\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Set File Attributes [C0050]\n  - File System::Get File Attributes [C0049]\n  - get file system object information\n  - Process::Create Thread [C0038]\n  - Process::Create Process [C0017]\n  - File System::Create File [C0016]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - set current directory\n  - find graphical window\n  - Discovery::File and Directory Discovery [E1083]\n  - link function at runtime on Windows\n  - File System::Copy File [C0045]\n  - File System::Move File [C0063]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Operating System::Registry::Query Registry Key [C0036.005]\n  - File System::Read File [C0051]\n  - File System::Delete File [C0047]\n\nHIGHLIGHTED CALLS:\n  GetTickCount\n  SetFileTime\n  GetTickCount64\n  GetAdaptersAddresses\n\n--- ARTIFACTS ---\n\nFILES DROPPED (5):\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHd.exe  [sha256: ab9b60de40fa8dda...]\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdLib.dll  [sha256: f08cecb6928548c5...]\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdLib64.dll  [sha256: 0c969c6d5b8f8bd8...]\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdHelper64.exe  [sha256: d859c9c92a460a3b...]\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\uninstall.exe  [sha256: a8efc0406dfba1ae...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHd.exe\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHdHelper64.exe\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHdLib.dll\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHdLib64.dll\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1200.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1210.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1220.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1230.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1240.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1250.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1260.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1270.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1280.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\Wallpapers\\1290.jpg\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\uninstall.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHd.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHdLib.dll\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHdLib64.dll\n  C:\\Users\\<USER>\\AppData\\Local\\Microsoft\\Windows\\Temporary Internet Files\\Content.IE5\\M9M2IBLU\\CheckNewVersion[1].txt\n  C:\\Users\\<USER>\\AppData\\Local\\Microsoft\\Windows\\Temporary Internet Files\\Content.IE5\\M9M2IBLU\\collect[1].gif\n\nPERSISTENCE REGISTRY KEYS (2):\n  Software\\Microsoft\\Windows\\CurrentVersion\\Run\\WallpaperHd = \"C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHd.exe\" /regrun\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\WallpaperHd = \"C:\\Documents and Settings\\Administrator\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHd.exe\n\nOTHER REGISTRY KEYS SET (10):\n  Software\\Microsoft\\Windows\\CurrentVersion\\Uninstall\\WallpaperHd\\DisplayName\n  Software\\Microsoft\\Windows\\CurrentVersion\\Uninstall\\WallpaperHd\\UninstallString\n  Software\\Microsoft\\Windows\\CurrentVersion\\Uninstall\\WallpaperHd\\Publisher\n  Software\\Microsoft\\Windows\\CurrentVersion\\Uninstall\\WallpaperHd\\DisplayVersion\n  Software\\Microsoft\\Windows\\CurrentVersion\\Uninstall\\WallpaperHd\\DisplayIcon\n  Control Panel\\Desktop\\LocalPersonalization\\Application\n  Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ProxyEnable\n  Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\Connections\\SavedLegacySettings\n  {6F0A4181-A97B-41A9-802B-9F04449E50C2}\\WpadDecisionReason\n  {6F0A4181-A97B-41A9-802B-9F04449E50C2}\\WpadDecisionTime\n\nMUTEXES CREATED (8):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ShimCacheMutex\n  oleacc-msaa-loaded\n\nPROCESSES CREATED (4):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHd.exe\" /install \"C:\\263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533\"\"\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\uninstall.exe\" /S _?=C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\"\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHd.exe /install 263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533.exe\n  C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdHelper64.exe /install\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA, GetProcAddress\n  ADVAPI32.dll: AdjustTokenPrivileges, OpenProcessToken\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Installer] Nullsoft Scriptable Install System\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'F-PROT': 'NSIS', 'Varist': 'NSIS'}\n\nPE SECTIONS (entropy):\n  .text        entropy=6.44\n  .rdata       entropy=5.24\n  .data        entropy=4.05\n  .ndata       entropy=0.0\n  .rsrc        entropy=4.57\n  [overlay]    entropy=7.996784687042236  size=7119630\n",
  "ground_truth": {
    "hash": "263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533",
    "verdict": "MALICIOUS",
    "detection_rate": 0.3553,
    "malicious_count": 27,
    "popular_threat_category": [
      {
        "value": "trojan",
        "count": 7
      },
      {
        "value": "adware",
        "count": 2
      }
    ],
    "popular_threat_name": [
      {
        "value": "wallpaperhd",
        "count": 6
      }
    ],
    "suggested_label": "trojan.wallpaperhd",
    "av_family_vote": "wallpaperhd",
    "has_family_label": true,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "query environment variable"
      },
      {
        "id": "T1083",
        "signature_description": "get file system object information"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry value"
      },
      {
        "id": "T1115",
        "signature_description": "open clipboard"
      },
      {
        "id": "T1125",
        "signature_description": "capture webcam image"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1134",
        "signature_description": "modify access privileges"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1497",
        "signature_description": "reference anti-VM strings targeting Xen"
      },
      {
        "id": "T1529",
        "signature_description": "shutdown system"
      },
      {
        "id": "T1547",
        "signature_description": "create shortcut via IShellLink"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "set file attributes",
      "get file attributes",
      "get file system object information",
      "create thread",
      "create process on Windows",
      "create or open file",
      "create or open registry key",
      "encode data using XOR",
      "contain loop",
      "query environment variable",
      "set current directory",
      "find graphical window",
      "enumerate files on Windows",
      "link function at runtime on Windows",
      "copy file",
      "move file",
      "query or enumerate registry value",
      "set registry value",
      "query or enumerate registry key",
      "read .ini file",
      "delete file",
      "create shortcut via IShellLink",
      "delete registry value",
      "delete registry key",
      "get file size",
      "check for time delay via GetTickCount",
      "set environment variable",
      "shutdown system",
      "get common file path",
      "accept command line arguments",
      "terminate process",
      "modify access privileges",
      "get OS version",
      "check OS version",
      "check if file exists",
      "get disk size",
      "capture webcam image",
      "open clipboard",
      "get graphical window text",
      "write clipboard data",
      "create directory",
      "delete directory",
      "enumerate files recursively",
      "read file on Windows",
      "write file on Windows",
      "hash data with CRC32",
      "reference anti-VM strings targeting Xen",
      "packaged as a NSIS installer",
      "(internal) installer file limitation"
    ],
    "iocs": {
      "ips": [
        "172.217.18.163",
        "216.58.206.14",
        "52.178.105.179",
        "67.26.131.254",
        "93.184.220.29"
      ],
      "domains": [
        "au-bg-shim.trafficmanager.net",
        "audownload.windowsupdate.nsatc.net",
        "auto.au.download.windowsupdate.com.c.footprint.net",
        "cs9.wac.phicdn.net",
        "ctldl.windowsupdate.com",
        "ocsp.digicert.com",
        "ocsp.pki.goog",
        "pki-goog.l.google.com",
        "wallhd.azurewebsites.net",
        "waws-prod-am2-123.cloudapp.net",
        "waws-prod-am2-123.vip.azurewebsites.windows.net",
        "www-google-analytics.l.google.com",
        "www.google-analytics.com"
      ],
      "urls": [
        "http://ctldl.windowsupdate.com/msdownload/update/v3/static/trustedr/en/disallowedcertstl.cab",
        "http://ctldl.windowsupdate.com/msdownload/update/v3/static/trustedr/en/disallowedcertstl.cab?53379364b6dd3cc4",
        "http://ocsp.digicert.com/MFEwTzBNMEswSTAJBgUrDgMCGgUABBTBL0V27RVZ7LBduom%2FnYB45SPUEwQU5Z1ZMIJHWMys%2BghUNoZ7OrUETfACEAiIzVJfGSRETRSlgpHeuVI%3D",
        "http://ocsp.pki.goog/gts1o1core/MFIwUDBOMEwwSjAJBgUrDgMCGgUABBRCRjDCJxnb3nDwj%2Fxz5aZfZjgXvAQUmNH4bhDrz5vsYJ8YkBug630J%2FSsCEQDM9j8%2B94pW6ggAAAAATZ8E",
        "https://wallhd.azurewebsites.net/Utils/CheckNewVersion?t=620&curVersion=1.4.0.12&target=newestSupported",
        "https://www.google-analytics.com/collect?v=1&tid=UA-96780277-2&cid=75c4d8f9-d671-46d9-86b1-77e51760d218&t=event&aip=1&ec=Wallpaper_HD_MA&ea=Install&el=20180701&z=1597667476"
      ],
      "registry_persistence": [
        "Software\\Microsoft\\Windows\\CurrentVersion\\Run\\WallpaperHd",
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\WallpaperHd"
      ],
      "registry_all_count": 21,
      "files_dropped": [
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHd.exe",
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdHelper64.exe",
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdLib.dll",
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdLib64.dll",
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\uninstall.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ShimCacheMutex",
        "oleacc-msaa-loaded"
      ],
      "processes": [
        "C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\WallpaperHd.exe\" /install \"C:\\263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533\"\"",
        "C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\\uninstall.exe\" /S _?=C:\\Documents and Settings\\<USER>\\Local Settings\\Application Data\\WallpaperHd\"",
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHd.exe /install 263499ff40d92efa30d45fc3ec41843f21d15cf6dee163618b7d06bc1b0aa533.exe",
        "C:\\Users\\<USER>\\AppData\\Local\\WallpaperHd\\WallpaperHdHelper64.exe /install"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox",
    "VirusTotal Observer"
  ]
}