# dataset/preprocessed/26638f7a5e84ff4b7d9521eecd18c549ab046a64007dfade7a94491c0ac54d60.json
{
  "hash": "26638f7a5e84ff4b7d9521eecd18c549ab046a64007dfade7a94491c0ac54d60",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  26638f7a5e84ff4b_beast catfight femdom .zip.exe\n    C:\\Users\\<USER>\\Downloads\\26638f7a5e84ff4b_beast catfight femdom .zip.exe\n    C:\\Users\\<USER>\\Downloads\\26638f7a5e84ff4b_beast catfight femdom .zip.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\26638f7a5e84ff4b7d9521eecd18c549ab046a64007dfade7a94491c0ac54d60(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\26638f7a5e84ff4b_beast catfight femdom .zip.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african beast full movie black hairunshaved .zip.exe  [sha256: 6d8d4aec75eb74d6...]\n  american animal hardcore girls balls .mpeg.exe  [sha256: 6fbd2653da702449...]\n  black animal sperm [free] ejaculation .mpeg.exe  [sha256: 58012ee1e3178e9e...]\n  black horse blowjob sleeping glans beautyfull .avi.exe  [sha256: 35ea9abceee4e987...]\n  blowjob full movie titts redhair  (Tatjana).mpg.exe  [sha256: cce0447e842ab6b1...]\n  brasilian kicking lesbian [milf] glans .mpg.exe  [sha256: fb54c7dcb246f6b7...]\n  bukkake full movie gorgeoushorny .zip.exe  [sha256: 83e63e2ec88878eb...]\n  danish beastiality lesbian hot (!) cock .mpeg.exe  [sha256: c5945169954b2518...]\n  danish cumshot horse sleeping glans shower .mpeg.exe  [sha256: 62290dca0c87ac85...]\n  danish gang bang sperm [bangbus] traffic .avi.exe  [sha256: 10d9b79b0238b82d...]\n  defrag.exe  [sha256: af7fb4b8ae94dde6...]\n  deinstall.exe  [sha256: 96c29085f92e7198...]\n  driver.exe  [sha256: 1e1b6eb9cb324047...]\n  drwatson32.exe  [sha256: 2f9e3e3387210a4c...]\n  freeftp.exe  [sha256: fd1ee8866902ba54...]\n  fucking licking femdom .rar.exe  [sha256: b1dc8ca7876dbbe4...]\n  fucking several models feet  (Jenna,Jade).mpeg.exe  [sha256: c488db996c07a558...]\n  fucking sleeping titts balls .mpg.exe  [sha256: eb2d983cffe42bf6...]\n  gay voyeur girly  (Ashley,Melissa).mpg.exe  [sha256: 10836ff22a73f726...]\n  hardcore public leather .mpg.exe  [sha256: 4bbfc587de94229f...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\swedish nude bukkake big  (Tatjana).rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\horse big cock .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\lingerie masturbation mistress .rar.exe\n  C:\\Documents and Settings\\All Users\\Templates\\black animal sperm [free] ejaculation .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\fucking licking femdom .rar.exe\n  C:\\Documents and Settings\\Default User\\Templates\\fucking sleeping titts balls .mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\brasilian kicking lesbian [milf] glans .mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\hardcore public leather .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\lingerie hidden glans hotel .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\sperm public glans fishy .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\italian beastiality fucking [free]  (Sylvia).avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\danish beastiality lesbian hot (!) cock .mpeg.exe\n  C:\\WINDOWS\\Temp\\russian beastiality lingerie full movie 50+ .avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\trambling public  (Janette).zip.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\russian fetish horse [milf] boots .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\bukkake full movie gorgeoushorny .zip.exe\n  C:\\WINDOWS\\assembly\\tmp\\black horse blowjob sleeping glans beautyfull .avi.exe\n  C:\\WINDOWS\\ime\\shared\\danish cumshot horse sleeping glans shower .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\gay voyeur girly  (Ashley,Melissa).mpg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.875821113586426  size=21606\n",
  "ground_truth": {
    "hash": "26638f7a5e84ff4b7d9521eecd18c549ab046a64007dfade7a94491c0ac54d60",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\danish beastiality xxx big cock blondie  (Jade).mpg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\indian fetish blowjob licking hole stockings  (Jade).mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\swedish porn fucking hidden hole shoes  (Samantha).rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\fucking hot (!) titts .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\brasilian horse fucking uncut 40+ .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\japanese kicking hardcore masturbation titts  (Anniston,Melissa).zip.exe",
        "C:\\Program Files\\dotnet\\shared\\american kicking beast public granny  (Gina,Tatjana).mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "african beast full movie black hairunshaved .zip.exe",
        "american animal hardcore girls balls .mpeg.exe",
        "black animal sperm [free] ejaculation .mpeg.exe",
        "black horse blowjob sleeping glans beautyfull .avi.exe",
        "blowjob full movie titts redhair  (Tatjana).mpg.exe",
        "brasilian kicking lesbian [milf] glans .mpg.exe",
        "bukkake full movie gorgeoushorny .zip.exe",
        "danish beastiality lesbian hot (!) cock .mpeg.exe",
        "danish cumshot horse sleeping glans shower .mpeg.exe",
        "danish gang bang sperm [bangbus] traffic .avi.exe",
        "defrag.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking licking femdom .rar.exe",
        "fucking several models feet  (Jenna,Jade).mpeg.exe",
        "fucking sleeping titts balls .mpg.exe",
        "gay voyeur girly  (Ashley,Melissa).mpg.exe",
        "hardcore public leather .mpg.exe",
        "horse big cock .mpeg.exe",
        "indian cumshot horse masturbation sweet .mpg.exe",
        "italian beastiality fucking [free]  (Sylvia).avi.exe",
        "italian nude sperm [milf] swallow  (Jenna,Melissa).zip.exe",
        "italian porn fucking catfight wifey  (Christine,Jade).zip.exe",
        "japanese fetish xxx public titts ash .rar.exe",
        "lingerie hidden glans hotel .zip.exe",
        "lingerie masturbation mistress .rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian beastiality lingerie full movie 50+ .avi.exe",
        "russian fetish horse [milf] boots .rar.exe",
        "screensaver.scr",
        "setup.exe",
        "soccer.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}