# dataset/preprocessed/2888d88c6d2e4e961c121844de612afde517e94898f8718663f0595b2f5b9cd7.json
{
  "hash": "2888d88c6d2e4e961c121844de612afde517e94898f8718663f0595b2f5b9cd7",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (7):\n  dns.msftncsi.com \u2192 131.107.255.255\n  ftp.mfrees.nokly.com \u2192 69.172.201.153\n  ftp.mfrees.nokly.com \u2192 76.223.54.146, 13.248.169.48\n  ftp.myfreesite.foreverhost.us\n  ftp.myfreesite.foreverhost.us \u2192 103.224.212.212\n  gripto.leadhoster.com\n  gripto.leadhoster.com \u2192 127.0.0.1\n\nIP TRAFFIC (2):\n  114.114.114.114:53 (UDP)\n  69.172.201.153:21 ()\n\nSERVICES OPENED (1):\n  RASMAN\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  <PATH_SAMPLE.EXE>\n    %ProgramFiles(x86)%\\6dfcc4e4\\jusched.exe\n  <SYSTEM32>\\svchost.exe\n  <SYSTEM32>\\svchost.exe\n  %ProgramFiles(x86)%\\6dfcc4e4\\jusched.exe\n\nCOMMAND EXECUTIONS (2):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - enumerate files on Windows\n  - parse PE header\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - Discovery::File and Directory Discovery [E1083]\n  - parse PE header\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - debug build\n  - contain a resource (.rsrc) section\n  - contains PDB path\n\n--- ARTIFACTS ---\n\nFILES DROPPED (11):\n  jusched.exe  [sha256: 845ad205110cdb55...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: 1699e1fd53b0b82b...]\n  %ProgramFiles(x86)%\\6dfcc4e4\\6dfcc4e4  [sha256: 795b20420c630ebd...]\n  %ProgramFiles(x86)%\\6dfcc4e4\\info_a  [sha256: 2785bbf60a3f89e4...]\n  %ProgramFiles(x86)%\\6dfcc4e4\\jusched.exe  [sha256: b4627362a69e5e39...]\n  <SYSTEM32>\\tasks\\update23  [sha256: 291e6671d07da1c9...]\n  %WINDIR%\\tasks\\update23.job  [sha256: 9002e992284f19ae...]\n  c:\\program files\\d460fff6\\jusched.exe  [sha256: b4f5ae4282b7911c...]\n  c:\\program files\\d460fff6\\d460fff6  [sha256: 976c041e9e2ad664...]\n  c:\\program files\\d460fff6\\info_a  [sha256: 9e160b3b15c572f2...]\n  c:\\windows\\tasks\\update23.job  [sha256: 38f1005362fc16cc...]\n\nFILES WRITTEN (7):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files (x86)\\71d2c87e\\info_a\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\info_a\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\WINDOWS\\Tasks\\Update23.job\n  C:\\Windows\\Tasks\\Update23.job\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:JavaUpdate23\n\nOTHER REGISTRY KEYS SET (10):\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\UNCAsIntranet\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\AutoDetect\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectId_\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\LruList\\CurrentLru\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectId\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectLru\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectLru_\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\ObjectTable\\E9\\_FileId_\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\ObjectTable\\E9\\_Usn_\n  \\REGISTRY\\A\\{c7d7b51e-fa86-11e9-a082-d6bfc30e5f1f}\\DefaultObjectStore\\ObjectTable\\E9\\_UsnJournalId_\n\nMUTEXES CREATED (13):\n  Local\\WininetConnectionMutex\n  Local\\WininetProxyRegistryMutex\n  Local\\WininetStartupMutex\n  Local\\ZoneAttributeCacheCounterMutex\n  Local\\ZonesCacheCounterMutex\n  Local\\ZonesCounterMutex\n  Local\\ZonesLockedCacheCounterMutex\n  Local\\_!MSFTHISTORY!_\n  Local\\c:!users!admin!appdata!local!microsoft!windows!history!history.ie5!\n  Local\\c:!users!admin!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n\nMUTEXES OPENED (9):\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{3f443105-1033-449b-9004-ede865c3013a}\n  \\Sessions\\1\\BaseNamedObjects\\Local\\WininetConnectionMutex\n  \\Sessions\\1\\BaseNamedObjects\\Local\\WininetProxyRegistryMutex\n  \\Sessions\\1\\BaseNamedObjects\\Local\\WininetStartupMutex\n  \\Sessions\\1\\BaseNamedObjects\\Local\\_!MSFTHISTORY!_\n  \\Sessions\\1\\BaseNamedObjects\\Local\\c:!users!user!appdata!local!microsoft!windows!history!history.ie5!\n  \\Sessions\\1\\BaseNamedObjects\\Local\\c:!users!user!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  \\Sessions\\1\\BaseNamedObjects\\Local\\c:!users!user!appdata!roaming!microsoft!windows!cookies!\n  \\Sessions\\1\\BaseNamedObjects\\RasPbFile\n\nPROCESSES CREATED (5):\n  %ProgramFiles(x86)%\\6dfcc4e4\\jusched.exe\n  <PATH_SAMPLE.EXE>\n  C:\\2888d88c6d2e4e961c121844de612afde517e94898f8718663f0595b2f5b9cd7.exe\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n  C:\\Program Files\\d460fff6\\jusched.exe \n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  [Overlay] PDB 2.0 file link\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.92\n  .rdata       entropy=4.65\n  .data        entropy=0.87\n  .idata       entropy=3.03\n  .rsrc        entropy=3.18\n  .reloc       entropy=4.24\n  [overlay]    entropy=7.051174163818359  size=7946\n",
  "ground_truth": {
    "hash": "2888d88c6d2e4e961c121844de612afde517e94898f8718663f0595b2f5b9cd7",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9306,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "enumerate files on Windows",
      "parse PE header",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "debug build",
      "contain a resource (.rsrc) section",
      "contains PDB path"
    ],
    "iocs": {
      "ips": [
        "114.114.114.114",
        "69.172.201.153"
      ],
      "domains": [
        "dns.msftncsi.com",
        "ftp.mfrees.nokly.com",
        "ftp.myfreesite.foreverhost.us",
        "gripto.leadhoster.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe"
      ],
      "registry_all_count": 52,
      "files_dropped": [
        "%ProgramFiles(x86)%\\6dfcc4e4\\6dfcc4e4",
        "%ProgramFiles(x86)%\\6dfcc4e4\\info_a",
        "%ProgramFiles(x86)%\\6dfcc4e4\\jusched.exe",
        "%WINDIR%\\tasks\\update23.job",
        "<SYSTEM32>\\tasks\\update23",
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "c:\\program files\\d460fff6\\d460fff6",
        "c:\\program files\\d460fff6\\info_a",
        "c:\\program files\\d460fff6\\jusched.exe",
        "c:\\windows\\tasks\\update23.job",
        "jusched.exe"
      ],
      "mutexes": [
        "Local\\WininetConnectionMutex",
        "Local\\WininetProxyRegistryMutex",
        "Local\\WininetStartupMutex",
        "Local\\ZoneAttributeCacheCounterMutex",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex",
        "Local\\_!MSFTHISTORY!_",
        "Local\\c:!users!admin!appdata!local!microsoft!windows!history!history.ie5!",
        "Local\\c:!users!admin!appdata!local!microsoft!windows!temporary internet files!content.ie5!",
        "Local\\c:!users!admin!appdata!roaming!microsoft!windows!cookies!",
        "RasPbFile",
        "ShimCacheMutex"
      ],
      "processes": [
        "%ProgramFiles(x86)%\\6dfcc4e4\\jusched.exe",
        "<PATH_SAMPLE.EXE>",
        "C:\\2888d88c6d2e4e961c121844de612afde517e94898f8718663f0595b2f5b9cd7.exe",
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"",
        "C:\\Program Files\\d460fff6\\jusched.exe "
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Dr.Web vxCube",
    "Sangfor ZSand",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}