# dataset/preprocessed/294148444521370fb4669573b80247553330473dd97cbad0cd736b1298bc431f.json
{
  "hash": "294148444521370fb4669573b80247553330473dd97cbad0cd736b1298bc431f",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  <PATH_SAMPLE.EXE>\n    %WINDIR%\\microsofthelp.exe\n  %WINDIR%\\microsofthelp.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (6):\n  - get common file path\n  - link function at runtime on Windows\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get MAC address on Windows\n  - get file size\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - reference absolute stream path on Windows\n  - Process::Create Process [C0017]\n  - File System::Create File [C0016]\n  - Discovery::File and Directory Discovery [E1083]\n  - contain loop\n  - File System::Delete File [C0047]\n  - Process::Terminate Process [C0018]\n  - Process::Create Mutex [C0042]\n  - Process::Check Mutex [C0043]\n  - link function at runtime on Windows\n  - Communication::HTTP Communication::Create Request [C0002.012]\n  - Communication::HTTP Communication::Open URL [C0002.004]\n  - Communication::HTTP Communication::Get Response [C0002.017]\n  - Command and Control::C2 Communication::Receive Data [B0030.002]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - get MAC address on Windows\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n\nHIGHLIGHTED CALLS:\n  {'file': ''}\n\n--- ARTIFACTS ---\n\nFILES DROPPED (3):\n  microsofthelp.exe  [sha256: 3c8c420586f45cb3...]\n  C:\\Windows\\microsofthelp.exe  [sha256: 14e750cf176bfd23...]\n  %WINDIR%\\microsofthelp.exe  [sha256: c5d02699f8a9b713...]\n\nFILES WRITTEN (4):\n  %WINDIR%\\hideplugin.dll\n  C:\\WINDOWS\\HidePlugin.dll\n  C:\\WINDOWS\\microsofthelp.exe\n  C:\\Windows\\HidePlugin.dll\n\nPERSISTENCE REGISTRY KEYS (2):\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\microsofthelp = %WINDIR%\\microsofthelp.exe\n  HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\microsofthelp = C:\\Windows\\microsofthelp.exe\n\nMUTEXES CREATED (3):\n  ShimCacheMutex\n  \\Sessions\\1\\BaseNamedObjects\\pomdfghrt\n  pomdfghrt\n\nMUTEXES OPENED (7):\n  \\Sessions\\1\\BaseNamedObjects\\Local\\WininetConnectionMutex\n  \\Sessions\\1\\BaseNamedObjects\\Local\\WininetProxyRegistryMutex\n  \\Sessions\\1\\BaseNamedObjects\\Local\\WininetStartupMutex\n  \\Sessions\\1\\BaseNamedObjects\\Local\\_!MSFTHISTORY!_\n  \\Sessions\\1\\BaseNamedObjects\\Local\\c:!users!user!appdata!local!microsoft!windows!history!history.ie5!\n  \\Sessions\\1\\BaseNamedObjects\\Local\\c:!users!user!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  \\Sessions\\1\\BaseNamedObjects\\Local\\c:!users!user!appdata!roaming!microsoft!windows!cookies!\n\nPROCESSES CREATED (33):\n  %WINDIR%\\microsofthelp.exe\n  <PATH_SAMPLE.EXE>\n  C:\\Users\\RA491~1.VUL\\AppData\\Local\\Temp\\acb7049b9bf5be1991491685c7fddd1b.exe \n  C:\\Users\\RA491~1.VUL\\AppData\\Local\\Temp\\llCt64.exe /Terminate /originalpath C:\\zjmbua\\bin\\llCt.exe\n  C:\\WINDOWS\\microsofthelp.exe\n  C:\\Windows\\System32\\slui.exe -Embedding\n  C:\\Windows\\System32\\wevtutil.exe query-events microsoft-windows-powershell/operational /rd:true /e:root /format:xml /uni:true\n  C:\\Windows\\microsofthelp.exe\n  C:\\Windows\\system32\\AUDIODG.EXE 0x684\n  C:\\Windows\\system32\\schtasks.exe /delete /f /TN Microsoft\\Windows\\Customer Experience Improvement Program\\Uploader\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA, GetProcAddress\n\nCOMPILER AND PACKER:\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'Cyren': 'dropped, dropped, dropped'}\n\nPE SECTIONS (entropy):\n  .text        entropy=6.4\n  .rdata       entropy=4.73\n  .data        entropy=0.96\n  .shoooo      entropy=7.84\n  .imports     entropy=3.61\n  [overlay]    entropy=4.34074592590332  size=19932\n",
  "ground_truth": {
    "hash": "294148444521370fb4669573b80247553330473dd97cbad0cd736b1298bc431f",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8451,
    "malicious_count": 60,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1082",
        "signature_description": "get MAC address on Windows"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "create thread",
      "create or open registry key",
      "reference absolute stream path on Windows",
      "create process on Windows",
      "create or open file",
      "get common file path",
      "contain loop",
      "delete file",
      "terminate process",
      "create mutex",
      "check mutex and exit",
      "link function at runtime on Windows",
      "create HTTP request",
      "connect to URL",
      "read data from Internet",
      "receive data",
      "set registry value",
      "query or enumerate registry value",
      "persist via Run registry key",
      "get MAC address on Windows",
      "read file on Windows",
      "write file on Windows",
      "get file size"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "<HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\microsofthelp",
        "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\microsofthelp"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "%WINDIR%\\microsofthelp.exe",
        "C:\\Windows\\microsofthelp.exe",
        "microsofthelp.exe"
      ],
      "mutexes": [
        "ShimCacheMutex",
        "\\Sessions\\1\\BaseNamedObjects\\pomdfghrt",
        "pomdfghrt"
      ],
      "processes": [
        "%WINDIR%\\microsofthelp.exe",
        "<PATH_SAMPLE.EXE>",
        "C:\\Users\\RA491~1.VUL\\AppData\\Local\\Temp\\acb7049b9bf5be1991491685c7fddd1b.exe ",
        "C:\\Users\\RA491~1.VUL\\AppData\\Local\\Temp\\llCt64.exe /Terminate /originalpath C:\\zjmbua\\bin\\llCt.exe",
        "C:\\WINDOWS\\microsofthelp.exe",
        "C:\\Windows\\System32\\slui.exe -Embedding",
        "C:\\Windows\\System32\\wevtutil.exe query-events microsoft-windows-powershell/operational /rd:true /e:root /format:xml /uni:true",
        "C:\\Windows\\microsofthelp.exe",
        "C:\\Windows\\system32\\AUDIODG.EXE 0x684",
        "C:\\Windows\\system32\\schtasks.exe /delete /f /TN Microsoft\\Windows\\Customer Experience Improvement Program\\Uploader",
        "C:\\Windows\\system32\\wbem\\wmiprvse.exe -Embedding",
        "C:\\zjmbua\\bin\\WhcylUNVp.exe inject 2692 2696 \\OMWMjpHGvBqGuiJOay.dll",
        "C:\\zjmbua\\bin\\WhcylUNVp.exe inject 2900 2904 \\OMWMjpHGvBqGuiJOay.dll",
        "C:\\zjmbua\\bin\\llCt.exe /Terminate",
        "\\??\\C:\\Windows\\system32\\conhost.exe -116456853489835296-942581350-3007093072087774494-2098449535-2084063021625693949",
        "\\??\\C:\\Windows\\system32\\conhost.exe -1837137579-958961245-470209709-280633883-147385235142171096-380352788309447548",
        "\\??\\C:\\Windows\\system32\\conhost.exe -20788411071260080611-18996266911272207592-72560066521032626511507669358-1606798987",
        "\\??\\C:\\Windows\\system32\\conhost.exe -390234220-14976914472556931071479241391609292196-612958801-1138572417156386931",
        "\\??\\C:\\Windows\\system32\\conhost.exe -620617709724892269-11605033241639922081-650309344-1997154181-3691361461000522025",
        "\\??\\C:\\Windows\\system32\\conhost.exe -782650330-1855995281772501949-776655312-168816281320788043621181829805-1396924481",
        "\\??\\C:\\Windows\\system32\\conhost.exe 1291864738107509361810285621514832475-11757483171027695934787580143-1586924108",
        "\\??\\C:\\Windows\\system32\\conhost.exe 1318694754-724123728343215771-51469153292603070815957903231489438145-288105000",
        "\\??\\C:\\Windows\\system32\\conhost.exe 1348420305240321265497584265-2025404794-9253534439965829429841762100875370",
        "\\??\\C:\\Windows\\system32\\conhost.exe 1574082304-16373545541491358244-52408538019940134496611446270173021180324278",
        "\\??\\C:\\Windows\\system32\\conhost.exe 19314939721843263201-453114434-1056079658-14788144172120867080-457347255-1481494035",
        "\\??\\C:\\Windows\\system32\\conhost.exe 26219948-1151537930-794992024-1589675703-1758972318-1906977180-612704063491314810",
        "\\??\\C:\\Windows\\system32\\conhost.exe 32207542-1957047445-13393695134775600431599467622122774810519111352591378865689",
        "\\??\\C:\\Windows\\system32\\conhost.exe 388891507-5681877522010371259-1911090075-17775235552986603191970290088-291567906",
        "\\??\\C:\\Windows\\system32\\conhost.exe 548870315-1854198121-1061881541-1084140135901061314-466655809-1988303558619858807",
        "\\??\\C:\\Windows\\system32\\conhost.exe 554991827-206682126647793316-1669306858-1961824023-590808710-1429026616-587894644",
        "\\??\\C:\\Windows\\system32\\conhost.exe 606530160118202842311754829171374100648-805116563-2123224122-1200309367938739665",
        "\\??\\C:\\Windows\\system32\\conhost.exe 9461853021071746445152040064314286611761312392353442638026-1384338111409218454",
        "wmiadap.exe /F /T /R"
      ]
    }
  },
  "num_sandboxes": 6,
  "sandbox_names": [
    "CAPA",
    "Dr.Web vxCube",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox",
    "Yomi Hunter"
  ]
}