# dataset/preprocessed/2a95af0c2dbdd5d20bcdd1d684d05b0da75e61e2677234ff8cd7f1c245a60c68.json
{
  "hash": "2a95af0c2dbdd5d20bcdd1d684d05b0da75e61e2677234ff8cd7f1c245a60c68",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  2a95af0c2dbdd5d2_tyrkish cum sperm girls .avi.exe\n    C:\\Users\\<USER>\\Downloads\\2a95af0c2dbdd5d2_tyrkish cum sperm girls .avi.exe\n    C:\\Users\\<USER>\\Downloads\\2a95af0c2dbdd5d2_tyrkish cum sperm girls .avi.exe\n      C:\\Users\\<USER>\\Downloads\\2a95af0c2dbdd5d2_tyrkish cum sperm girls .avi.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\2a95af0c2dbdd5d20bcdd1d684d05b0da75e61e2677234ff8cd7f1c245a60c68(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\2a95af0c2dbdd5d2_tyrkish cum sperm girls .avi.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american animal blowjob girls glans .rar.exe  [sha256: 3f540af6fc8d27b4...]\n  american fetish horse public titts black hairunshaved  (Sylvia).rar.exe  [sha256: 739c098e8ba92340...]\n  beast full movie ejaculation .avi.exe  [sha256: b822e882d7e52c75...]\n  beast girls cock .rar.exe  [sha256: cb1b0669d281e4ce...]\n  beast voyeur .mpeg.exe  [sha256: be1aa96f94dbbd81...]\n  beast [bangbus] feet bondage .mpg.exe  [sha256: d15fb9fd8e226a5d...]\n  black nude blowjob [milf] hole .zip.exe  [sha256: 21be0721dd37683c...]\n  black porn horse sleeping .avi.exe  [sha256: 9b536d3730580f8e...]\n  brasilian nude blowjob voyeur  (Melissa).mpg.exe  [sha256: c91cd0f567dede9b...]\n  british trambling full movie .zip.exe  [sha256: 1afe3e065586b385...]\n  bukkake voyeur  (Jade).rar.exe  [sha256: 739cdc83ede96886...]\n  danish cum gay catfight glans girly .avi.exe  [sha256: 38037e7f6ae77ee5...]\n  defrag.exe  [sha256: ce776388243913c2...]\n  defrag64.exe  [sha256: d29aa4c8add7a2c5...]\n  driver.exe  [sha256: 45edc0effc9d5488...]\n  drwatson32.exe  [sha256: 0d32ccafa463e681...]\n  fucking full movie  (Curtney).zip.exe  [sha256: 494100ab78d1700e...]\n  hardcore uncut feet stockings .zip.exe  [sha256: b397e1730f46b228...]\n  hardcore uncut titts  (Kathrin,Samantha).rar.exe  [sha256: 9201e46192af0f16...]\n  IE9setup.exe  [sha256: 43dc73c066d41bcd...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\tyrkish beastiality trambling lesbian hole gorgeoushorny .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\beast [bangbus] feet bondage .mpg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\american animal blowjob girls glans .rar.exe\n  C:\\Documents and Settings\\All Users\\Templates\\lesbian [bangbus] .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\hardcore uncut feet stockings .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\fucking full movie  (Curtney).zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\american fetish horse public titts black hairunshaved  (Sylvia).rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\lesbian uncut girly .rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\beast voyeur .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\hardcore uncut titts  (Kathrin,Samantha).rar.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\british trambling full movie .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\bukkake voyeur  (Jade).rar.exe\n  C:\\WINDOWS\\Temp\\italian nude beast [milf] .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\russian beastiality blowjob [milf] redhair .rar.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\black nude blowjob [milf] hole .zip.exe\n  C:\\WINDOWS\\assembly\\temp\\beast full movie ejaculation .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\danish cum gay catfight glans girly .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\black porn horse sleeping .avi.exe\n  C:\\WINDOWS\\security\\templates\\russian beastiality trambling catfight fishy  (Anniston,Sylvia).rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880164623260498  size=1976301\n",
  "ground_truth": {
    "hash": "2a95af0c2dbdd5d20bcdd1d684d05b0da75e61e2677234ff8cd7f1c245a60c68",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9437,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\xxx uncut  (Samantha).zip.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\gay uncut .zip.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\black action trambling [milf] cock .mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\japanese handjob blowjob big .rar.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\bukkake girls stockings .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\xxx [bangbus] high heels .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\american cumshot fucking uncut .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\black fetish horse several models gorgeoushorny  (Ashley,Sylvia).rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\fucking licking  (Tatjana).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\lesbian [free] hole granny .zip.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\russian horse hardcore several models titts shower .zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\japanese porn lingerie licking traffic  (Britney,Sylvia).avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\danish porn trambling [milf] cock .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\lingerie voyeur titts .mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\lingerie girls swallow .avi.exe",
        "C:\\Program Files\\dotnet\\shared\\japanese action horse hidden .rar.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american animal blowjob girls glans .rar.exe",
        "american fetish horse public titts black hairunshaved  (Sylvia).rar.exe",
        "beast [bangbus] feet bondage .mpg.exe",
        "beast full movie ejaculation .avi.exe",
        "beast girls cock .rar.exe",
        "beast voyeur .mpeg.exe",
        "black nude blowjob [milf] hole .zip.exe",
        "black porn horse sleeping .avi.exe",
        "brasilian nude blowjob voyeur  (Melissa).mpg.exe",
        "british trambling full movie .zip.exe",
        "bukkake voyeur  (Jade).rar.exe",
        "danish cum gay catfight glans girly .avi.exe",
        "defrag.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking full movie  (Curtney).zip.exe",
        "hardcore uncut feet stockings .zip.exe",
        "hardcore uncut titts  (Kathrin,Samantha).rar.exe",
        "indian cumshot bukkake hot (!) feet shoes .mpg.exe",
        "indian gang bang fucking [bangbus] ejaculation  (Anniston,Melissa).mpg.exe",
        "italian nude beast [milf] .rar.exe",
        "japanese nude horse big redhair .mpg.exe",
        "lesbian [bangbus] .zip.exe",
        "lesbian uncut girly .rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "rundll64.exe",
        "russian beastiality blowjob [milf] redhair .rar.exe",
        "russian beastiality trambling catfight fishy  (Anniston,Sylvia).rar.exe",
        "screensaver.scr",
        "sperm voyeur feet latex .rar.exe",
        "start.exe",
        "swedish handjob blowjob masturbation hole blondie .avi.exe",
        "tyrkish beastiality trambling lesbian hole gorgeoushorny .avi.exe",
        "undelete.exe",
        "xxx public glans .mpeg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}