# dataset/preprocessed/2b337550dbcc817d4d71433cb2179657d99f90045d1ce0b6bcdac7e38f82faee.json
{
  "hash": "2b337550dbcc817d4d71433cb2179657d99f90045d1ce0b6bcdac7e38f82faee",
  "report_text": "--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (27):\n  - encode data using XOR\n  - hide graphical window\n  - set file attributes\n  - access PEB ldr_data\n  - check Internet connectivity via WinINet\n  - accept command line arguments\n  - query or enumerate registry value\n  - get geographical location\n  - get file size\n  - find graphical window\n  - get disk size\n  - link function at runtime on Windows\n  - get system information on Windows\n  - capture screenshot\n  - log keystrokes via polling\n  - link many functions at runtime\n  - log keystrokes\n  - get keyboard layout\n  - enumerate gui resources\n  - query or enumerate registry key\n  - get common file path\n  - get local IPv4 addresses\n  - get hostname\n  - enumerate network shares\n  - persist via Run registry key\n  - get disk information\n  - reference Base64 string\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Process::Create Thread [C0038]\n  - File System::Create File [C0016]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - hide graphical window\n  - set global application hook\n  - inspect load icon resource\n  - File System::Set File Attributes [C0050]\n  - reference absolute stream path on Windows\n  - access PEB ldr_data\n  - Process::Open Process [C0065]\n  - check Internet connectivity via WinINet\n  - contain loop\n  - Execution::Command and Scripting Interpreter [E1059]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - get geographical location\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  Age of Sail 2 Patch.exe  [sha256: 34d4eedd3550322c...]\n  AIM Account Stealer Crack.exe  [sha256: e516f7af017292bb...]\n  AIM Account Stealer Patch.exe  [sha256: 586d5aaf993bf5e8...]\n  Battle.net Patch.exe  [sha256: b91686ebfccee85b...]\n  CKY3 - Bam Margera World Industries Alien Workshop Key Generator.exe  [sha256: 61c3901888962cb5...]\n  CloneCD Full Downloader.exe  [sha256: 5f1cb05908ad1176...]\n  Combat Flight Simulator 3 Key Generator.exe  [sha256: 5879ad8af7ceb077...]\n  Crazy Taxi Key Generator.exe  [sha256: f5ea42d34f4a646a...]\n  Critical Point Manga game Patch.exe  [sha256: 9d0569badb248957...]\n  Deadly Dozen Crack.exe  [sha256: 19eedb1d3e4ccd04...]\n  Dweebs 2 Key Generator.exe  [sha256: f6ddd212c467693c...]\n  Emperor Rise Of the Middle Kingdom Full Downloader.exe  [sha256: 84bb038662da5d6a...]\n  F1 Grand Pix 4 ISO - Full Downloader.exe  [sha256: 245e18b3bb8dbaee...]\n  F1 Grand Pix 4 Patch.exe  [sha256: e121733df0187da4...]\n  Freedom Force ISO - Full Downloader.exe  [sha256: 3ac88c0801e2a8b6...]\n  Grand Theft Auto 3 Patch.exe  [sha256: d4afae035e6f6034...]\n  GTA 3 Full Downloader.exe  [sha256: 7dca93c3d7e4c96e...]\n  GTA3 Full Downloader.exe  [sha256: 481f7f4074ad8bea...]\n  Half-life ONLINE ISO - Full Downloader.exe  [sha256: 69f328ed72b33851...]\n  Industry Giant 2 Crack.exe  [sha256: 6eff37309cebd537...]\n\nFILES WRITTEN (30):\n  C:\\2b337550dbcc817d4d71433cb2179657d99f90045d1ce0b6bcdac7e38f82faee\n  C:\\My Downloads\\AIM Account Stealer Crack.exe\n  C:\\My Downloads\\AIM Account Stealer Patch.exe\n  C:\\My Downloads\\Battle.net Patch.exe\n  C:\\My Downloads\\CKY3 - Bam Margera World Industries Alien Workshop Key Generator.exe\n  C:\\My Downloads\\Combat Flight Simulator 3 Key Generator.exe\n  C:\\My Downloads\\Critical Point Manga game Patch.exe\n  C:\\My Downloads\\GTA 3 Full Downloader.exe\n  C:\\My Downloads\\GTA3 Full Downloader.exe\n  C:\\My Downloads\\Half-life ONLINE ISO - Full Downloader.exe\n  C:\\My Downloads\\Industry Giant 2 Crack.exe\n  C:\\My Downloads\\Macromedia Crack.exe\n  C:\\My Downloads\\Macromedia Dreamweaver MX Crack.exe\n  C:\\My Downloads\\Macromedia Full Downloader.exe\n  C:\\My Downloads\\Macromedia Key Generator.exe\n  C:\\My Downloads\\Mafia Crack.exe\n  C:\\My Downloads\\Necromania Trap Of Darkness Full Downloader.exe\n  C:\\My Downloads\\Neverwinter Nights Full Downloader.exe\n  C:\\My Downloads\\Norton Utilities 2002 XP ISO - Full Downloader.exe\n  C:\\My Downloads\\Soldiers Of Anarchy Crack.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\VIR-HELP = C:\\Windows\\KRrq2XY.exe\n\nMUTEXES CREATED (1):\n  I-Worm.PlutonX\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n\nPE SECTIONS (entropy):\n  CODE         entropy=6.5\n  DATA         entropy=4.83\n  BSS          entropy=0.0\n  .idata       entropy=4.97\n  .tls         entropy=0.0\n  .rdata       entropy=0.2\n  .reloc       entropy=0.0\n  .rsrc        entropy=3.97\n  [overlay]    entropy=7.989641189575195  size=44309\n",
  "ground_truth": {
    "hash": "2b337550dbcc817d4d71433cb2179657d99f90045d1ce0b6bcdac7e38f82faee",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9014,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1016",
        "signature_description": "check Internet connectivity via WinINet"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via polling"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "generate random numbers using the Delphi LCG",
      "create or open registry key",
      "create thread",
      "create or open file",
      "delay execution",
      "PEB access",
      "encode data using XOR",
      "allocate RWX memory",
      "hide graphical window",
      "set global application hook",
      "inspect load icon resource",
      "set file attributes",
      "reference absolute stream path on Windows",
      "access PEB ldr_data",
      "open process",
      "check Internet connectivity via WinINet",
      "contain loop",
      "accept command line arguments",
      "query or enumerate registry value",
      "get geographical location",
      "read file on Windows",
      "write file on Windows",
      "get file size",
      "set thread local storage value",
      "get thread local storage value",
      "find graphical window",
      "delete file",
      "get disk size",
      "create directory",
      "get OS version",
      "link function at runtime on Windows",
      "extract resource via kernel32 functions",
      "get system information on Windows",
      "capture screenshot",
      "log keystrokes via polling",
      "link many functions at runtime",
      "log keystrokes",
      "get keyboard layout",
      "set application hook",
      "enumerate gui resources",
      "get graphical window text",
      "read .ini file",
      "query or enumerate registry key",
      "set registry value",
      "copy file",
      "get common file path",
      "receive data on socket",
      "receive data",
      "get local IPv4 addresses",
      "resolve DNS",
      "get hostname",
      "terminate process",
      "enumerate network shares",
      "persist via Run registry key",
      "get disk information",
      "create mutex",
      "check mutex and exit",
      "contain a thread local storage (.tls) section",
      "contain a resource (.rsrc) section",
      "reference Base64 string",
      "compiled with Borland Delphi",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\VIR-HELP"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "AIM Account Stealer Crack.exe",
        "AIM Account Stealer Patch.exe",
        "Age of Sail 2 Patch.exe",
        "Battle.net Patch.exe",
        "CKY3 - Bam Margera World Industries Alien Workshop Key Generator.exe",
        "CloneCD Full Downloader.exe",
        "Combat Flight Simulator 3 Key Generator.exe",
        "Crazy Taxi Key Generator.exe",
        "Critical Point Manga game Patch.exe",
        "Deadly Dozen Crack.exe",
        "Dweebs 2 Key Generator.exe",
        "Emperor Rise Of the Middle Kingdom Full Downloader.exe",
        "F1 Grand Pix 4 ISO - Full Downloader.exe",
        "F1 Grand Pix 4 Patch.exe",
        "Freedom Force ISO - Full Downloader.exe",
        "GTA 3 Full Downloader.exe",
        "GTA3 Full Downloader.exe",
        "Grand Theft Auto 3 Patch.exe",
        "Half-life ONLINE ISO - Full Downloader.exe",
        "Industry Giant 2 Crack.exe",
        "International Cricket Captain 2003 Patch.exe",
        "KaZaA Spyware Remover ISO - Full Downloader.exe",
        "LordOfTheRingsr Key Generator.exe",
        "Macromedia Crack.exe",
        "Macromedia Dreamweaver MX Crack.exe",
        "Macromedia Full Downloader.exe",
        "Macromedia Key Generator.exe",
        "Mafia Crack.exe",
        "Necromania Trap Of Darkness Full Downloader.exe",
        "Neverwinter Nights Full Downloader.exe",
        "Neverwinter Nights Key Generator.exe",
        "Norton AntiVirus 2002 Key Generator.exe",
        "Norton Utilities 2002 XP ISO - Full Downloader.exe",
        "Shakira Patch.exe",
        "Soldier Of Fortune 2 Crack.exe",
        "Soldiers Of Anarchy Crack.exe",
        "Star Wars Starfighter Patch.exe",
        "Strike Fighter Project 1 ISO - Full Downloader.exe",
        "Sudden Strike 2 Patch.exe",
        "The Eye Of Kraken Crack.exe",
        "The Eye Of Kraken ISO - Full Downloader.exe",
        "The Eye Of Kraken Patch.exe",
        "Unreal Tournament 3 Patch.exe",
        "Warcraft 3 ONLINE Patch.exe",
        "Warcraft 3 battle.net Patch.exe",
        "Windows XP SP1 Patch.exe",
        "Winzip 8.0 ISO - Full Downloader.exe",
        "Winzip 8.0 Patch.exe",
        "Xbox.info ISO - Full Downloader.exe",
        "ZoneAlarm Firewall Patch.exe"
      ],
      "mutexes": [
        "I-Worm.PlutonX"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}