# dataset/preprocessed/2c119126adf484a82872528c22c597f3011079f089ca8603d4c48250ed8753d7.json
{
  "hash": "2c119126adf484a82872528c22c597f3011079f089ca8603d4c48250ed8753d7",
  "report_text": "--- NETWORK ACTIVITY ---\n\nSERVICES OPENED (1):\n  Messenger\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  2c119126adf484a8_swedish nude masturbation feet .rar.exe\n    C:\\Users\\<USER>\\Downloads\\2c119126adf484a8_swedish nude masturbation feet .rar.exe\n\nCOMMAND EXECUTIONS (1):\n  C:\\Users\\<USER>\\Downloads\\2c119126adf484a8_swedish nude masturbation feet .rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (9):\n  C:\\Windows\\mssrv.exe  [sha256: aee630455621c298...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\trambling hidden nipples girly  (Curtney).mpeg.exe  [sha256: 6a0cd1b10506a665...]\n  C:\\Program Files\\dotnet\\shared\\black animal gay [milf] traffic .mpg.exe  [sha256: b55192379c3448ed...]\n  C:\\Program Files\\DVD Maker\\Shared\\tyrkish hardcore gay masturbation .avi.exe  [sha256: abad6f1435c76fd1...]\n  C:\\Program Files\\Windows Journal\\Templates\\hardcore lingerie lesbian hairy .rar.exe  [sha256: ed061b1152b05a83...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\animal [milf] glans .rar.exe  [sha256: 279a4329ebaeaa32...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\british gay big sm  (Liz).rar.exe  [sha256: f69d84d908baf76e...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\nude cum uncut vagina  (Kathrin).mpeg.exe  [sha256: 376e6264c7f732cb...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\asian hardcore action uncut hairy .zip.exe  [sha256: 11235cfc2f186604...]\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880829811096191  size=1016851\n",
  "ground_truth": {
    "hash": "2c119126adf484a82872528c22c597f3011079f089ca8603d4c48250ed8753d7",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9437,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\british gay big sm  (Liz).rar.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\nude cum uncut vagina  (Kathrin).mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\asian hardcore action uncut hairy .zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\trambling hidden nipples girly  (Curtney).mpeg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\tyrkish hardcore gay masturbation .avi.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\hardcore lingerie lesbian hairy .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\animal [milf] glans .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\black animal gay [milf] traffic .mpg.exe",
        "C:\\Windows\\mssrv.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Jujubox"
  ]
}