# dataset/preprocessed/2c5e08da0424069c7ad1d612dfeefedda0e5fe72062d5686322d4a5fa93fad66.json
{
  "hash": "2c5e08da0424069c7ad1d612dfeefedda0e5fe72062d5686322d4a5fa93fad66",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  2c5e08da0424069c_japanese cum sperm licking blondie .zip.exe\n    C:\\Users\\<USER>\\Downloads\\2c5e08da0424069c_japanese cum sperm licking blondie .zip.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\2c5e08da0424069c7ad1d612dfeefedda0e5fe72062d5686322d4a5fa93fad66(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\2c5e08da0424069c_japanese cum sperm licking blondie .zip.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african bukkake kicking [bangbus]  (Sandy).mpeg.exe  [sha256: 78cb43e2aaf3966a...]\n  african cum gang bang [milf] .avi.exe  [sha256: b35f6010aab23375...]\n  american gay kicking catfight 50+ .zip.exe  [sha256: f46f1d5a79ff16bd...]\n  american xxx licking vagina  (Sonja).avi.exe  [sha256: c232d91acdec1ac8...]\n  animal girls .rar.exe  [sha256: 389decb847686887...]\n  beast full movie  (Sandy,Britney).rar.exe  [sha256: 6e0e9c6d761d8e33...]\n  black blowjob lesbian uncut .rar.exe  [sha256: 063c7c544c573861...]\n  cum uncut .avi.exe  [sha256: 256a87dcc48dd3ae...]\n  danish beast big legs  (Britney,Gina).zip.exe  [sha256: 2df683a0c9a63d07...]\n  defrag.exe  [sha256: b26c9c35932cc17c...]\n  deinstall.exe  [sha256: eddcea231e10f260...]\n  edit.exe  [sha256: 79d50be75b0d4342...]\n  fetish blowjob [bangbus] balls .zip.exe  [sha256: 7acca9be118b9088...]\n  freeftp.exe  [sha256: c7a10031abf87b60...]\n  french nude full movie legs  (Sandy,Ashley).avi.exe  [sha256: e2d1357497a51d09...]\n  fucking girls titts  (Gina,Jenna).rar.exe  [sha256: 53767e855956da1a...]\n  gay hot (!) ash femdom .mpeg.exe  [sha256: d352716364e022e8...]\n  handjob fetish voyeur shoes  (Liz).rar.exe  [sha256: 008f5cb694dc2a04...]\n  hdcleaner.exe  [sha256: 91504c96a217f033...]\n  IE9setup.exe  [sha256: 5a4d9c9f7ef39b42...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\tyrkish porn fetish masturbation .zip.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\fetish blowjob [bangbus] balls .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\danish beast big legs  (Britney,Gina).zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\spanish action voyeur lady .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\cum uncut .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\swedish handjob full movie boobs ejaculation .mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\french nude full movie legs  (Sandy,Ashley).avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\italian handjob handjob voyeur cock traffic  (Janette,Kathrin).avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\norwegian fucking sperm masturbation nipples castration .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\swedish handjob trambling hidden .mpg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\swedish cumshot public .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\russian beast catfight nipples traffic  (Tatjana,Kathrin).zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\gay hot (!) ash femdom .mpeg.exe\n  C:\\WINDOWS\\Temp\\animal girls .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\swedish bukkake cumshot lesbian hole gorgeoushorny .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\fucking girls titts  (Gina,Jenna).rar.exe\n  C:\\WINDOWS\\assembly\\temp\\kicking sperm [milf] hole .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\porn sperm several models redhair .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\trambling several models .mpg.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.880416393280029  size=1272534\n",
  "ground_truth": {
    "hash": "2c5e08da0424069c7ad1d612dfeefedda0e5fe72062d5686322d4a5fa93fad66",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "compiled with Go",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\malaysia sperm catfight .mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\spanish xxx beastiality catfight .rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\blowjob lesbian several models vagina hotel .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\danish gay beast several models boobs  (Karin,Jade).mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\american trambling voyeur .mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\horse hot (!) .mpg.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "african bukkake kicking [bangbus]  (Sandy).mpeg.exe",
        "african cum gang bang [milf] .avi.exe",
        "american gay kicking catfight 50+ .zip.exe",
        "american xxx licking vagina  (Sonja).avi.exe",
        "animal girls .rar.exe",
        "beast full movie  (Sandy,Britney).rar.exe",
        "black blowjob lesbian uncut .rar.exe",
        "cum uncut .avi.exe",
        "danish beast big legs  (Britney,Gina).zip.exe",
        "defrag.exe",
        "deinstall.exe",
        "edit.exe",
        "fetish blowjob [bangbus] balls .zip.exe",
        "freeftp.exe",
        "french nude full movie legs  (Sandy,Ashley).avi.exe",
        "fucking girls titts  (Gina,Jenna).rar.exe",
        "gay hot (!) ash femdom .mpeg.exe",
        "handjob fetish voyeur shoes  (Liz).rar.exe",
        "hdcleaner.exe",
        "install32.exe",
        "italian handjob handjob voyeur cock traffic  (Janette,Kathrin).avi.exe",
        "kicking sperm [milf] hole .mpeg.exe",
        "lesbian [free] .mpeg.exe",
        "mssrv.exe",
        "norwegian fucking sperm masturbation nipples castration .mpeg.exe",
        "notes.txt.exe",
        "porn sperm several models redhair .mpeg.exe",
        "readme.txt.exe",
        "russian beast catfight nipples traffic  (Tatjana,Kathrin).zip.exe",
        "setup.exe",
        "soccer.exe",
        "spanish action voyeur lady .zip.exe",
        "start.exe",
        "summergames.exe",
        "swedish bukkake cumshot lesbian hole gorgeoushorny .mpg.exe",
        "swedish cumshot public .zip.exe",
        "swedish fetish [bangbus] .rar.exe",
        "swedish handjob full movie boobs ejaculation .mpg.exe",
        "swedish handjob trambling hidden .mpg.exe",
        "trambling kicking lesbian  (Anniston,Sonja).rar.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}