# dataset/preprocessed/2cc1dcb88baa2cc05f8268f66e15d5fd266bb99f9d85c23587f9422678d8dc2d.json
{
  "hash": "2cc1dcb88baa2cc05f8268f66e15d5fd266bb99f9d85c23587f9422678d8dc2d",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\2cc1dcb88baa2cc05f8268f66e15d5fd266bb99f9d85c23587f9422678d8dc2d(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\2cc1dcb88baa2cc0_trambling xxx public boobs 40+  (anniston).zip.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  african fucking blowjob sleeping titts bondage .avi.exe  [sha256: 6fa5a36658f387cc...]\n  african handjob fetish [free]  (Samantha).avi.exe  [sha256: 3416aff54a801cd5...]\n  animal masturbation hairy .avi.exe  [sha256: 0c52f12783eec1f1...]\n  asian fetish animal [free] nipples sm  (Janette).mpeg.exe  [sha256: 69b333b1b1d62680...]\n  beast beastiality voyeur high heels .avi.exe  [sha256: a9834a429d14a724...]\n  black horse masturbation ash .mpg.exe  [sha256: f6fad28945154952...]\n  brasilian horse horse uncut mistress .mpeg.exe  [sha256: e7c2b78100f4325c...]\n  brasilian lingerie handjob voyeur femdom  (Curtney).avi.exe  [sha256: c331fbaccc58a0e6...]\n  brasilian sperm licking fishy .avi.exe  [sha256: 01116470401b8d64...]\n  british horse sleeping femdom .zip.exe  [sha256: 97a64302066a9624...]\n  bukkake hardcore lesbian .mpg.exe  [sha256: 949523787d66ab88...]\n  canadian lingerie xxx sleeping circumcision .avi.exe  [sha256: 3f25ad51c353e062...]\n  cum girls ash .mpeg.exe  [sha256: cb68068d1eafda82...]\n  cum handjob hot (!) vagina lady .zip.exe  [sha256: 18f6e39a5f40b318...]\n  defrag64.exe  [sha256: 004f2122202fbed8...]\n  freeftp.exe  [sha256: 69a116e0eb8f75e8...]\n  german cum [free]  (Liz).zip.exe  [sha256: c58d598639a25c22...]\n  hdcleaner.exe  [sha256: 0736401b9928331d...]\n  horse hot (!) .zip.exe  [sha256: 4496c2e6007307a3...]\n  install.exe  [sha256: 01d6db0c88ca2ac2...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\tyrkish cumshot xxx sleeping lady .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\beast beastiality voyeur high heels .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\swedish horse full movie ash circumcision .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\norwegian xxx lesbian boots .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\african handjob fetish [free]  (Samantha).avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\brasilian horse horse uncut mistress .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\british horse sleeping femdom .zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\cum handjob hot (!) vagina lady .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\trambling lesbian girly .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\asian fetish animal [free] nipples sm  (Janette).mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\malaysia nude beast hot (!) boots .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\porn full movie .mpg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\norwegian nude catfight shower  (Liz,Sonja).mpeg.exe\n  C:\\WINDOWS\\Temp\\german cum [free]  (Liz).zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\nude [free] .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\swedish kicking xxx [bangbus] titts hairy .avi.exe\n  C:\\WINDOWS\\assembly\\temp\\cum girls ash .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\brasilian lingerie handjob voyeur femdom  (Curtney).avi.exe\n  C:\\WINDOWS\\ime\\shared\\african fucking blowjob sleeping titts bondage .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.879925727844238  size=2051977\n",
  "ground_truth": {
    "hash": "2cc1dcb88baa2cc05f8268f66e15d5fd266bb99f9d85c23587f9422678d8dc2d",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9412,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\american kicking blowjob voyeur  (Curtney).mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\russian animal sperm masturbation cock  (Jenna,Jade).mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\indian horse horse catfight feet castration .zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\sperm sleeping lady .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\indian porn hardcore hot (!) hole ejaculation .zip.exe",
        "C:\\Program Files\\dotnet\\shared\\swedish animal bukkake several models cock femdom .rar.exe",
        "C:\\Windows\\mssrv.exe",
        "african fucking blowjob sleeping titts bondage .avi.exe",
        "african handjob fetish [free]  (Samantha).avi.exe",
        "animal masturbation hairy .avi.exe",
        "asian fetish animal [free] nipples sm  (Janette).mpeg.exe",
        "beast beastiality voyeur high heels .avi.exe",
        "black horse masturbation ash .mpg.exe",
        "brasilian horse horse uncut mistress .mpeg.exe",
        "brasilian lingerie handjob voyeur femdom  (Curtney).avi.exe",
        "brasilian sperm licking fishy .avi.exe",
        "british horse sleeping femdom .zip.exe",
        "bukkake hardcore lesbian .mpg.exe",
        "canadian lingerie xxx sleeping circumcision .avi.exe",
        "cum girls ash .mpeg.exe",
        "cum handjob hot (!) vagina lady .zip.exe",
        "defrag64.exe",
        "freeftp.exe",
        "german cum [free]  (Liz).zip.exe",
        "hdcleaner.exe",
        "horse hot (!) .zip.exe",
        "install.exe",
        "lesbian porn uncut mature  (Janette).mpg.exe",
        "malaysia nude beast hot (!) boots .avi.exe",
        "mssrv.exe",
        "norwegian nude catfight shower  (Liz,Sonja).mpeg.exe",
        "norwegian trambling trambling uncut ash ejaculation .rar.exe",
        "norwegian xxx lesbian boots .mpeg.exe",
        "nude [free] .mpg.exe",
        "pgp9.exe",
        "porn full movie .mpg.exe",
        "rundll64.exe",
        "russian animal several models ash swallow .avi.exe",
        "screensaver.scr",
        "setup.exe",
        "spanish cum beastiality voyeur mature .avi.exe",
        "start.exe",
        "summergames.exe",
        "swedish horse full movie ash circumcision .zip.exe",
        "swedish kicking xxx [bangbus] titts hairy .avi.exe",
        "trambling lesbian girly .mpeg.exe",
        "tyrkish cumshot xxx sleeping lady .mpg.exe",
        "unformat.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}