# dataset/preprocessed/2cf4e366479132f17d411d707b996e7d34b857bbbec08829a0012a2296d4b313.json
{
  "hash": "2cf4e366479132f17d411d707b996e7d34b857bbbec08829a0012a2296d4b313",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\2cf4e366479132f17d411d707b996e7d34b857bbbec08829a0012a2296d4b313(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\2cf4e366479132f1_mssrv.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american fetish blowjob catfight titts upskirt .rar.exe  [sha256: c0dc6ab3dc3b6a5e...]\n  american handjob fucking girls boots  (Jenna,Janette).mpeg.exe  [sha256: 4282831172b1bd21...]\n  american kicking lesbian hot (!) ash .mpg.exe  [sha256: 4dfd216b60f3c023...]\n  beast hot (!) bondage .avi.exe  [sha256: 3ca92df9c2305c1b...]\n  brasilian beastiality fucking girls glans mistress .mpeg.exe  [sha256: 47c8c858460f6c12...]\n  brasilian cumshot fucking hidden hole boots  (Jade).rar.exe  [sha256: fdfef047939d1a06...]\n  brasilian handjob sperm full movie high heels .mpeg.exe  [sha256: ae4a58db0026289e...]\n  british bukkake hot (!) .zip.exe  [sha256: 0f9a45717bedd79b...]\n  bukkake sleeping .avi.exe  [sha256: 1e50532ee69ce8f0...]\n  danish animal sperm [bangbus] glans gorgeoushorny .rar.exe  [sha256: d8f7ec92aab30497...]\n  danish beastiality bukkake public glans .avi.exe  [sha256: 378d9c0b44ebe480...]\n  defrag.exe  [sha256: a66e1367bd9e8c5d...]\n  defrag64.exe  [sha256: 2da34cf2c7ddeb5f...]\n  driver.exe  [sha256: 23c366fc90b52bd1...]\n  drwatson32.exe  [sha256: 5a3a5f1b3067c722...]\n  freeftp.exe  [sha256: 5d9c7895cf0901db...]\n  fucking sleeping hole .mpg.exe  [sha256: ef6773eb85f4461b...]\n  fucking voyeur  (Tatjana).avi.exe  [sha256: 163e6cd50cf6debf...]\n  gay voyeur wifey .rar.exe  [sha256: d9a8c81ddf21c899...]\n  horse hidden mature .zip.exe  [sha256: 18551ba2191fb615...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\danish beastiality bukkake public glans .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\italian porn lingerie several models .rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\gay voyeur wifey .rar.exe\n  C:\\Documents and Settings\\All Users\\Templates\\swedish cum bukkake public hotel .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\xxx hidden lady  (Ashley,Sylvia).rar.exe\n  C:\\Documents and Settings\\Default User\\Templates\\danish animal sperm [bangbus] glans gorgeoushorny .rar.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\fucking voyeur  (Tatjana).avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\horse hidden mature .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\japanese horse lingerie [free]  (Jade).mpg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\bukkake sleeping .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\italian cum xxx [milf] traffic .mpg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\xxx masturbation glans .rar.exe\n  C:\\WINDOWS\\Temp\\japanese animal gay sleeping granny .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\russian gang bang horse sleeping .zip.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\british bukkake hot (!) .zip.exe\n  C:\\WINDOWS\\assembly\\temp\\beast hot (!) bondage .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\brasilian beastiality fucking girls glans mistress .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\trambling hidden swallow  (Britney,Sylvia).zip.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\american fetish blowjob catfight titts upskirt .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880168437957764  size=1983094\n",
  "ground_truth": {
    "hash": "2cf4e366479132f17d411d707b996e7d34b857bbbec08829a0012a2296d4b313",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9275,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\gay kicking hot (!) vagina girly  (Sylvia,Anniston).zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\malaysia bukkake trambling masturbation hole .mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\danish trambling licking beautyfull .rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\fucking beastiality lesbian 40+ .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\british beast horse [milf] cock circumcision .zip.exe",
        "C:\\Program Files\\dotnet\\shared\\danish beast nude sleeping redhair  (Melissa).rar.exe",
        "C:\\Windows\\mssrv.exe",
        "american fetish blowjob catfight titts upskirt .rar.exe",
        "american handjob fucking girls boots  (Jenna,Janette).mpeg.exe",
        "american kicking lesbian hot (!) ash .mpg.exe",
        "beast hot (!) bondage .avi.exe",
        "brasilian beastiality fucking girls glans mistress .mpeg.exe",
        "brasilian cumshot fucking hidden hole boots  (Jade).rar.exe",
        "brasilian handjob sperm full movie high heels .mpeg.exe",
        "british bukkake hot (!) .zip.exe",
        "bukkake sleeping .avi.exe",
        "danish animal sperm [bangbus] glans gorgeoushorny .rar.exe",
        "danish beastiality bukkake public glans .avi.exe",
        "defrag.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking sleeping hole .mpg.exe",
        "fucking voyeur  (Tatjana).avi.exe",
        "gay voyeur wifey .rar.exe",
        "horse hidden mature .zip.exe",
        "italian cum xxx [milf] traffic .mpg.exe",
        "italian fetish horse [bangbus] cock redhair .rar.exe",
        "italian porn lingerie several models .rar.exe",
        "japanese animal gay sleeping granny .rar.exe",
        "japanese beastiality bukkake catfight .zip.exe",
        "japanese horse lingerie [free]  (Jade).mpg.exe",
        "mssrv.exe",
        "rundll64.exe",
        "russian gang bang horse sleeping .zip.exe",
        "screensaver.scr",
        "setup.exe",
        "soccer.exe",
        "swedish cum bukkake public hotel .mpeg.exe",
        "swedish porn gay licking .mpg.exe",
        "trambling hidden swallow  (Britney,Sylvia).zip.exe",
        "tyrkish beastiality bukkake lesbian hotel .avi.exe",
        "xxx hidden lady  (Ashley,Sylvia).rar.exe",
        "xxx masturbation glans .rar.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}