# dataset/preprocessed/2ed81ac3d0ffaad63703f0af5c65fc99e4d639d3304da4ca34aa380b08f9d227.json
{
  "hash": "2ed81ac3d0ffaad63703f0af5c65fc99e4d639d3304da4ca34aa380b08f9d227",
  "report_text": "--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (28):\n  - encode data using XOR\n  - hide graphical window\n  - set file attributes\n  - access PEB ldr_data\n  - check Internet connectivity via WinINet\n  - accept command line arguments\n  - query or enumerate registry value\n  - get geographical location\n  - get file size\n  - find graphical window\n  - get disk size\n  - link function at runtime on Windows\n  - get system information on Windows\n  - capture screenshot\n  - log keystrokes via polling\n  - link many functions at runtime\n  - log keystrokes\n  - get keyboard layout\n  - enumerate gui resources\n  - query or enumerate registry key\n  - get common file path\n  - get local IPv4 addresses\n  - get hostname\n  - enumerate network shares\n  - persist via Run registry key\n  - get disk information\n  - packed with UPX\n  - reference Base64 string\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Cryptography::Generate Pseudo-random Sequence [C0021]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Process::Create Thread [C0038]\n  - File System::Create File [C0016]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - hide graphical window\n  - set global application hook\n  - inspect load icon resource\n  - File System::Set File Attributes [C0050]\n  - reference absolute stream path on Windows\n  - access PEB ldr_data\n  - Process::Open Process [C0065]\n  - check Internet connectivity via WinINet\n  - contain loop\n  - Execution::Command and Scripting Interpreter [E1059]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - get geographical location\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n\n--- ARTIFACTS ---\n\nFILES WRITTEN (30):\n  C:\\My Downloads\\AIM Account Stealer ISO - Full Downloader.exe\n  C:\\My Downloads\\Age of Sail 2 Key Generator.exe\n  C:\\My Downloads\\Austerlitz Napoleons Greatest Victory Patch.exe\n  C:\\My Downloads\\BORLAND Delphi 7 Crack.exe\n  C:\\My Downloads\\Battle.net Patch.exe\n  C:\\My Downloads\\Cabelas Ultimate Deer Hunt 2 Crack.exe\n  C:\\My Downloads\\Clive Barker\u62af Undying Patch.exe\n  C:\\My Downloads\\Dark Age Of Camelot Shrouded Isles ISO - Full Downloader.exe\n  C:\\My Downloads\\Dark Age Of Camelot Shrouded Isles Key Generator.exe\n  C:\\My Downloads\\Dark Age Of Camelot Shrouded Isles Patch.exe\n  C:\\My Downloads\\Dweebs 2 Patch.exe\n  C:\\My Downloads\\Elder Scrolls III Morrowind THX Brrbrr Key Generator.exe\n  C:\\My Downloads\\Empire Earth Patch.exe\n  C:\\My Downloads\\GTA 3 Full Downloader.exe\n  C:\\My Downloads\\Hitman 2 Silent Assassin Crack.exe\n  C:\\My Downloads\\International Cricket Captain 2003 ISO - Full Downloader.exe\n  C:\\My Downloads\\MoviezChannelsInstaler Key Generator.exe\n  C:\\My Downloads\\Nero Burning Rom 5.8.0.1 Full Downloader.exe\n  C:\\My Downloads\\Norton AntiVirus 2002 Key Generator.exe\n  C:\\My Downloads\\Norton Utilities 2002 XP Key Generator.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\NMAIN = C:\\WINDOWS\\LaPx6csWFUBjFC.exe\u0000\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows NT\\CurrentVersion\\Windows\\Run\n\nMUTEXES CREATED (9):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  I-Worm.PlutonX\n  MSCTF.Shared.MUTEX.AEH\n  MSCTF.Shared.MUTEX.EBH\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n\nPE SECTIONS (entropy):\n  UPX0         entropy=6.54\n  UPX1         entropy=4.59\n  .rsrc        entropy=3.44\n  .imports     entropy=4.65\n  [overlay]    entropy=7.989109039306641  size=36982\n",
  "ground_truth": {
    "hash": "2ed81ac3d0ffaad63703f0af5c65fc99e4d639d3304da4ca34aa380b08f9d227",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8732,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1016",
        "signature_description": "check Internet connectivity via WinINet"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via polling"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk size"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "allocate RW memory",
      "generate random numbers using the Delphi LCG",
      "create or open registry key",
      "create thread",
      "create or open file",
      "delay execution",
      "PEB access",
      "encode data using XOR",
      "allocate RWX memory",
      "hide graphical window",
      "set global application hook",
      "inspect load icon resource",
      "set file attributes",
      "reference absolute stream path on Windows",
      "access PEB ldr_data",
      "open process",
      "check Internet connectivity via WinINet",
      "contain loop",
      "accept command line arguments",
      "query or enumerate registry value",
      "get geographical location",
      "read file on Windows",
      "write file on Windows",
      "get file size",
      "set thread local storage value",
      "get thread local storage value",
      "find graphical window",
      "delete file",
      "get disk size",
      "create directory",
      "get OS version",
      "link function at runtime on Windows",
      "extract resource via kernel32 functions",
      "get system information on Windows",
      "capture screenshot",
      "log keystrokes via polling",
      "link many functions at runtime",
      "log keystrokes",
      "get keyboard layout",
      "set application hook",
      "enumerate gui resources",
      "get graphical window text",
      "read .ini file",
      "query or enumerate registry key",
      "set registry value",
      "copy file",
      "get common file path",
      "receive data on socket",
      "receive data",
      "get local IPv4 addresses",
      "resolve DNS",
      "get hostname",
      "terminate process",
      "enumerate network shares",
      "persist via Run registry key",
      "get disk information",
      "create mutex",
      "check mutex and exit",
      "packed with UPX",
      "contain a resource (.rsrc) section",
      "reference Base64 string",
      "compiled with Borland Delphi",
      "(internal) packer file limitation",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\NMAIN"
      ],
      "registry_all_count": 2,
      "files_dropped": [],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "I-Worm.PlutonX",
        "MSCTF.Shared.MUTEX.AEH",
        "MSCTF.Shared.MUTEX.EBH"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 2,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO"
  ]
}