# dataset/preprocessed/2fd359a3975176a29c64ba5a7125d48df8eca6130a18475374c9e2c85f3d6e9a.json
{
  "hash": "2fd359a3975176a29c64ba5a7125d48df8eca6130a18475374c9e2c85f3d6e9a",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  331d57be861d2ba73745485cb2a2e9e2.virus.exe\n    C:\\Windows\\CTS.exe\n\nCOMMAND EXECUTIONS (2):\n  C:\\Users\\Olivia\\AppData\\Local\\Temp\\N3NQt8Ll17DOwK6.exe\n  C:\\WINDOWS\\CTS.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (13):\n  - encode data using XOR\n  - access PEB ldr_data\n  - encrypt data using RC4 PRGA\n  - get file size\n  - query environment variable\n  - persist via Run registry key\n  - enumerate files on Windows\n  - enumerate files recursively\n  - get disk information\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - packed with UPX\n  - reference Base64 string\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - File System::Create File [C0016]\n  - Process::Create Process [C0017]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - access PEB ldr_data\n  - Process::Create Thread [C0038]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - Data::Checksum::CRC32 [C0032.001]\n  - Cryptography::Generate Pseudo-random Sequence::Use API [C0021.003]\n  - Cryptography::Encrypt Data::RC4 [C0027.009]\n  - Cryptography::Generate Pseudo-random Sequence::RC4 PRGA [C0021.004]\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Read File [C0051]\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Process::Terminate Process [C0018]\n  - File System::Writes File [C0052]\n  - Execution::Install Additional Program [B0023]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Process::Create Mutex [C0042]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (3):\n  422gWszzLjms4Ry.exe  [sha256: 6dc5020a170867ea...]\n  CTS.exe  [sha256: 27697ac78b171a47...]\n  C:\\Windows\\CTS.exe  [sha256: 27697ac78b171a47...]\n\nFILES WRITTEN (14):\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\422gWszzLjms4Ry.exe\n  C:\\Documents and Settings\\<USER>\\My Documents\\progs\\spotify.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Macromedia\\Flash Player\\www.macromedia.com\\bin\\fpupdateax\\fpupdateax.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouExplorer\\Temp\\sogouexplorerup.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\SogouPY\\SogouExplorer.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\Download\\P2PSetup.exe.tdl\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\Download\\QQPhotoDrawExSetupForQQ.exe.tdl\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\BackupDLTmp\\Download\\QzoneMusicInstall.exe.tdl\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\STemp\\SetupEx~0\\QQSetupEx.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\SafeBase\\QQSafeUD.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Tencent\\QQ\\commonf_inst\\TXSSOSetup.exe\n  C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\f9mDxWhZ6nlCcoR.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\rnJJEXU4B7Rdqgk.exe\n  C:\\WINDOWS\\CTS.exe\n\nPERSISTENCE REGISTRY KEYS (3):\n  HKLM\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\RUN = CTS\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\CTS = C\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\CTS = C:\\WINDOWS\\CTS.exe\u0000\n\nMUTEXES CREATED (2):\n  Global\\3pc6RWOgectGTFqCowxjeGy3XIGPtLwNrsr2zDctYD4hAU5pj4GW7rm8gHrHyTB6\n  ShimCacheMutex\n\nMUTEXES OPENED (1):\n  Local\\ShimViewer\n\nPROCESSES CREATED (3):\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\422gWszzLjms4Ry.exe\n  C:\\Users\\Olivia\\AppData\\Local\\Temp\\N3NQt8Ll17DOwK6.exe\n  C:\\WINDOWS\\CTS.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: CreateProcessW, GetProcAddress, LoadLibraryExW, LoadLibraryW\n\nCOMPILER AND PACKER:\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n\nPE SECTIONS (entropy):\n  UPX0         entropy=5.52\n  UPX1         entropy=3.12\n  UPX2         entropy=2.9\n  .imports     entropy=4.05\n  .reloc       entropy=6.53\n  [overlay]    entropy=4.168295383453369  size=24\n",
  "ground_truth": {
    "hash": "2fd359a3975176a29c64ba5a7125d48df8eca6130a18475374c9e2c85f3d6e9a",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8889,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "agent",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1082",
        "signature_description": "query environment variable"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "encode data using XOR",
      "create or open file",
      "create process on Windows",
      "create or open registry key",
      "PEB access",
      "access PEB ldr_data",
      "create thread",
      "delay execution",
      "contain loop",
      "hash data with CRC32",
      "generate random numbers via WinAPI",
      "encrypt data using RC4 PRGA",
      "get file size",
      "read file on Windows",
      "query environment variable",
      "set registry value",
      "terminate process",
      "write file on Windows",
      "write and execute a file",
      "persist via Run registry key",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get disk information",
      "create mutex",
      "check mutex and exit",
      "get thread local storage value",
      "print debug messages",
      "link function at runtime on Windows",
      "link many functions at runtime",
      "packed with UPX",
      "reference Base64 string",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\RUN",
        "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\CTS",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\CTS"
      ],
      "registry_all_count": 3,
      "files_dropped": [
        "422gWszzLjms4Ry.exe",
        "C:\\Windows\\CTS.exe",
        "CTS.exe"
      ],
      "mutexes": [
        "Global\\3pc6RWOgectGTFqCowxjeGy3XIGPtLwNrsr2zDctYD4hAU5pj4GW7rm8gHrHyTB6",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\422gWszzLjms4Ry.exe",
        "C:\\Users\\Olivia\\AppData\\Local\\Temp\\N3NQt8Ll17DOwK6.exe",
        "C:\\WINDOWS\\CTS.exe"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Lastline",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}