# dataset/preprocessed/2fd3e17f3fa64e23e009363be033ad911d500e917ba641f76505d6cd5fbe7695.json
{
  "hash": "2fd3e17f3fa64e23e009363be033ad911d500e917ba641f76505d6cd5fbe7695",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  2fd3e17f3fa64e23_danish bukkake several models ejaculation .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\2fd3e17f3fa64e23_danish bukkake several models ejaculation .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\2fd3e17f3fa64e23_danish bukkake several models ejaculation .mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\2fd3e17f3fa64e23e009363be033ad911d500e917ba641f76505d6cd5fbe7695(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\2fd3e17f3fa64e23_danish bukkake several models ejaculation .mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american sperm sleeping  (Liz).zip.exe  [sha256: fe861025975ec44c...]\n  animal action [milf] glans .mpeg.exe  [sha256: dbe214b80d6112bc...]\n  asian trambling porn hidden  (Anniston,Sylvia).mpeg.exe  [sha256: 44ba5f445520b199...]\n  beastiality uncut titts .avi.exe  [sha256: f26de3988ea4e660...]\n  black sperm voyeur titts circumcision .avi.exe  [sha256: a773b33a807623ca...]\n  blowjob sleeping ash .rar.exe  [sha256: 5735a399001cf940...]\n  brasilian handjob masturbation sweet  (Melissa).rar.exe  [sha256: 6fe95923f95d606c...]\n  cum trambling [free] granny .rar.exe  [sha256: a861fc8e2f800022...]\n  cumshot beast big femdom .rar.exe  [sha256: c5a524977feaec74...]\n  cumshot beast voyeur circumcision .zip.exe  [sha256: 32d951ba1f75ecf6...]\n  cumshot full movie  (Sandy).rar.exe  [sha256: 5d2a6f706388e5da...]\n  cumshot horse full movie wifey  (Melissa,Anniston).rar.exe  [sha256: 1068916316ec3908...]\n  cumshot horse sleeping titts latex .zip.exe  [sha256: 6570cd42a48e19a9...]\n  cumshot trambling catfight beautyfull  (Janette,Jenna).zip.exe  [sha256: d7522dd4875a39e7...]\n  cumshot trambling [free] .mpg.exe  [sha256: fb8f4883cd63ce02...]\n  defrag.exe  [sha256: fc2724f58a4ee6a1...]\n  driver.exe  [sha256: 914cee20fc714ec9...]\n  drwatson32.exe  [sha256: 9ba0dab24a82ea4b...]\n  edit.exe  [sha256: 0a933bcf8bf6cd58...]\n  explorer.exe  [sha256: 2d37ed4614b0d950...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\cumshot trambling [free] .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\asian trambling porn hidden  (Anniston,Sylvia).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\russian cum big nipples boots  (Kathrin).zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\gang bang lesbian sleeping  (Liz).mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\blowjob sleeping ash .rar.exe\n  C:\\Documents and Settings\\Default User\\Templates\\black sperm voyeur titts circumcision .avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\beastiality uncut titts .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\cumshot horse full movie wifey  (Melissa,Anniston).rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\cum trambling [free] granny .rar.exe\n  C:\\Program Files\\Movie Maker\\Shared\\fucking catfight .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\fetish bukkake [free] 40+ .avi.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\american sperm sleeping  (Liz).zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\gang bang horse [free] shower .zip.exe\n  C:\\WINDOWS\\Temp\\cumshot horse sleeping titts latex .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\italian xxx nude uncut girly  (Christine,Liz).mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\kicking fetish [free] lady .zip.exe\n  C:\\WINDOWS\\assembly\\temp\\sperm big .mpg.exe\n  C:\\WINDOWS\\assembly\\tmp\\xxx hot (!) shower .avi.exe\n  C:\\WINDOWS\\ime\\shared\\lingerie licking legs 40+ .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.880363464355469  size=725966\n",
  "ground_truth": {
    "hash": "2fd3e17f3fa64e23e009363be033ad911d500e917ba641f76505d6cd5fbe7695",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\french cumshot [bangbus] castration .mpeg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\danish sperm voyeur  (Sarah,Melissa).mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\lingerie full movie boots .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\british gay big castration .mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\asian beastiality hidden feet redhair  (Curtney).rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\indian lesbian public circumcision .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\malaysia fetish beast public boots .zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\british horse animal [bangbus] ash  (Janette,Gina).rar.exe",
        "C:\\Program Files\\dotnet\\shared\\fetish cum [free] cock .rar.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american sperm sleeping  (Liz).zip.exe",
        "animal action [milf] glans .mpeg.exe",
        "asian trambling porn hidden  (Anniston,Sylvia).mpeg.exe",
        "beastiality uncut titts .avi.exe",
        "black sperm voyeur titts circumcision .avi.exe",
        "blowjob sleeping ash .rar.exe",
        "brasilian handjob masturbation sweet  (Melissa).rar.exe",
        "cum trambling [free] granny .rar.exe",
        "cumshot beast big femdom .rar.exe",
        "cumshot beast voyeur circumcision .zip.exe",
        "cumshot full movie  (Sandy).rar.exe",
        "cumshot horse full movie wifey  (Melissa,Anniston).rar.exe",
        "cumshot horse sleeping titts latex .zip.exe",
        "cumshot trambling [free] .mpg.exe",
        "cumshot trambling catfight beautyfull  (Janette,Jenna).zip.exe",
        "defrag.exe",
        "driver.exe",
        "drwatson32.exe",
        "edit.exe",
        "explorer.exe",
        "fetish bukkake [free] 40+ .avi.exe",
        "fetish kicking hidden  (Ashley,Jade).rar.exe",
        "freeftp.exe",
        "fucking catfight .rar.exe",
        "fucking catfight 40+  (Curtney).avi.exe",
        "game3d.exe",
        "gang bang horse [free] shower .zip.exe",
        "gang bang lesbian sleeping  (Liz).mpeg.exe",
        "italian xxx nude uncut girly  (Christine,Liz).mpg.exe",
        "kicking fetish [free] lady .zip.exe",
        "lingerie licking legs 40+ .rar.exe",
        "mssrv.exe",
        "pgp9.exe",
        "russian cum big nipples boots  (Kathrin).zip.exe",
        "screensaver.scr",
        "soccer.exe",
        "speedup.exe",
        "sperm big .mpg.exe",
        "swedish trambling xxx [milf] stockings .mpeg.exe",
        "uninstall.exe",
        "word.exe",
        "xxx hot (!) shower .avi.exe",
        "xxx sleeping black hairunshaved .mpeg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}