# dataset/preprocessed/32332242ef02f3e12ccaa4b6800ff5c59a525b754c3a87582c9937780b94cd30.json
{
  "hash": "32332242ef02f3e12ccaa4b6800ff5c59a525b754c3a87582c9937780b94cd30",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  32332242ef02f3e1_horse public hole sm .avi.exe\n    C:\\Users\\<USER>\\Downloads\\32332242ef02f3e1_horse public hole sm .avi.exe\n    C:\\Users\\<USER>\\Downloads\\32332242ef02f3e1_horse public hole sm .avi.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\32332242ef02f3e12ccaa4b6800ff5c59a525b754c3a87582c9937780b94cd30(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\32332242ef02f3e1_horse public hole sm .avi.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (22):\n  action blowjob hot (!) hole .mpeg.exe  [sha256: 0e6f205343a5c78c...]\n  action [milf]  (Jade,Samantha).rar.exe  [sha256: 6587baa3c0919adb...]\n  animal [milf] .mpg.exe  [sha256: 69dfc963f3eebc1d...]\n  beastiality beast catfight  (Liz,Samantha).avi.exe  [sha256: 53d0969e4cef30c1...]\n  blowjob cumshot [bangbus] .avi.exe  [sha256: f68f7b0bee0aad36...]\n  brasilian bukkake fetish hidden traffic .zip.exe  [sha256: db7179bbd0ec1d74...]\n  british xxx big glans .rar.exe  [sha256: b4cce5ab5a0224dd...]\n  bukkake public wifey .avi.exe  [sha256: fce6a8ecd9a534fd...]\n  chinese gang bang masturbation  (Curtney).zip.exe  [sha256: 88059687f986f10b...]\n  chinese handjob horse [bangbus] vagina  (Kathrin).mpeg.exe  [sha256: fede9c475c8b4941...]\n  cum sleeping titts fishy .zip.exe  [sha256: f3ad36fd12df2415...]\n  danish xxx masturbation .avi.exe  [sha256: 6a1664d58a25ad58...]\n  defrag.exe  [sha256: 5e6ada597bc72b7f...]\n  C:\\Windows\\mssrv.exe  [sha256: b3e36462a337b8cf...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\spanish animal uncut sm .mpg.exe  [sha256: 78fcad0abdde9ac7...]\n  C:\\Program Files\\dotnet\\shared\\action beastiality [bangbus] feet castration .rar.exe  [sha256: e0dac81596d270c3...]\n  C:\\Program Files\\Windows Journal\\Templates\\african bukkake girls pregnant  (Britney).mpg.exe  [sha256: 8760223f4fc647b9...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\black fucking several models nipples .zip.exe  [sha256: 595415dbac230d9b...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\handjob nude big castration  (Christine,Jade).mpg.exe  [sha256: 11390e29781c6d18...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\british fetish full movie feet bedroom  (Kathrin).rar.exe  [sha256: f270fcaea7081d59...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\blowjob cumshot [bangbus] .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\norwegian blowjob masturbation boots .rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\russian gay beast sleeping glans hotel  (Karin).mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\japanese porn gay full movie .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\malaysia cumshot lesbian big nipples 50+ .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\fucking voyeur .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\action blowjob hot (!) hole .mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\malaysia lingerie gang bang full movie hotel .mpeg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\danish xxx masturbation .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\handjob voyeur feet balls .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\chinese handjob horse [bangbus] vagina  (Kathrin).mpeg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\cum sleeping titts fishy .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\norwegian trambling [milf] bedroom .rar.exe\n  C:\\WINDOWS\\Temp\\malaysia gay gay catfight boobs upskirt .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\lesbian licking .rar.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\beastiality beast catfight  (Liz,Samantha).avi.exe\n  C:\\WINDOWS\\assembly\\temp\\italian gay [free] boots  (Britney).mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\british xxx big glans .rar.exe\n  C:\\WINDOWS\\ime\\shared\\brasilian bukkake fetish hidden traffic .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.879981517791748  size=1192263\n",
  "ground_truth": {
    "hash": "32332242ef02f3e12ccaa4b6800ff5c59a525b754c3a87582c9937780b94cd30",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\handjob nude big castration  (Christine,Jade).mpg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\british fetish full movie feet bedroom  (Kathrin).rar.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\russian fucking kicking masturbation vagina .avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\spanish cumshot cumshot [bangbus] vagina  (Britney).avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\spanish animal uncut sm .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\african bukkake girls pregnant  (Britney).mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\black fucking several models nipples .zip.exe",
        "C:\\Program Files\\dotnet\\shared\\action beastiality [bangbus] feet castration .rar.exe",
        "C:\\Windows\\mssrv.exe",
        "action [milf]  (Jade,Samantha).rar.exe",
        "action blowjob hot (!) hole .mpeg.exe",
        "animal [milf] .mpg.exe",
        "beastiality beast catfight  (Liz,Samantha).avi.exe",
        "blowjob cumshot [bangbus] .avi.exe",
        "brasilian bukkake fetish hidden traffic .zip.exe",
        "british xxx big glans .rar.exe",
        "bukkake public wifey .avi.exe",
        "chinese gang bang masturbation  (Curtney).zip.exe",
        "chinese handjob horse [bangbus] vagina  (Kathrin).mpeg.exe",
        "cum sleeping titts fishy .zip.exe",
        "danish xxx masturbation .avi.exe",
        "defrag.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}