# dataset/preprocessed/3265da7a003709975d418e82ab7f9f4042676cc22a6074e814c97c4c3d9831c3.json
{
  "hash": "3265da7a003709975d418e82ab7f9f4042676cc22a6074e814c97c4c3d9831c3",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  lesbian several models mistress .rar.exe\n    C:\\Users\\<USER>\\Downloads\\lesbian several models mistress .rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\3265da7a003709975d418e82ab7f9f4042676cc22a6074e814c97c4c3d9831c3(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\lesbian several models mistress .rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  beast catfight titts traffic .rar.exe  [sha256: 53ce965c0cd0d267...]\n  beast hidden hairy  (Sonja,Curtney).avi.exe  [sha256: ca57fc5c1c075921...]\n  beast masturbation .mpeg.exe  [sha256: 3b5ba04a93a083fe...]\n  beast public titts .mpeg.exe  [sha256: b962b395bacdb3d7...]\n  black fetish fucking [milf]  (Jade).avi.exe  [sha256: 2f86dab3ddf993b6...]\n  black fetish lingerie [bangbus] titts circumcision .mpg.exe  [sha256: 6f5bdf70db634a3c...]\n  black porn sperm [free] .mpg.exe  [sha256: 1e555d5db96b2ff4...]\n  blowjob hot (!) black hairunshaved .mpeg.exe  [sha256: 92e55d977d1384a5...]\n  brasilian cumshot bukkake [bangbus] titts gorgeoushorny  (Jade).mpeg.exe  [sha256: 4f7f91367ca3a766...]\n  brasilian gang bang hardcore [bangbus] cock .zip.exe  [sha256: b9e9b67337a1bd29...]\n  danish cumshot hardcore uncut  (Sylvia).zip.exe  [sha256: ff6f972264cd149a...]\n  danish handjob blowjob masturbation stockings .mpeg.exe  [sha256: f3e4947bb136d354...]\n  danish kicking lesbian big pregnant .mpeg.exe  [sha256: a5b7bb71751bbd9d...]\n  danish nude lingerie [bangbus] leather .mpg.exe  [sha256: 76df23a7ef3204bf...]\n  defrag.exe  [sha256: 4d122f1640103d11...]\n  defrag64.exe  [sha256: 5e14e1171cbfea6c...]\n  deinstall.exe  [sha256: caebf626b1f45e6d...]\n  driver.exe  [sha256: 85efbe0aeeafeea0...]\n  drwatson32.exe  [sha256: a40e3b0e66d55dab...]\n  fucking lesbian redhair .mpg.exe  [sha256: 559c08ee83b6d0c8...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\beast public titts .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\horse [free] feet mistress  (Melissa).rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\black fetish fucking [milf]  (Jade).avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\swedish animal lingerie masturbation .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\tyrkish cum lingerie public feet castration .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\brasilian gang bang hardcore [bangbus] cock .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\swedish cum gay uncut titts .mpeg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\brasilian cumshot bukkake [bangbus] titts gorgeoushorny  (Jade).mpeg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\fucking lesbian redhair .mpg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\beast catfight titts traffic .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\swedish horse sperm big girly  (Sonja,Janette).mpeg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\italian cum trambling hidden leather  (Sandy,Sylvia).mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\danish kicking lesbian big pregnant .mpeg.exe\n  C:\\WINDOWS\\Temp\\kicking xxx hot (!) hotel .mpeg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\lesbian hot (!) gorgeoushorny  (Ashley,Jade).avi.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\black fetish lingerie [bangbus] titts circumcision .mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\xxx uncut .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\indian action xxx [milf] circumcision .zip.exe\n  C:\\WINDOWS\\ime\\shared\\gay several models young .avi.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880442142486572  size=1882774\n",
  "ground_truth": {
    "hash": "3265da7a003709975d418e82ab7f9f4042676cc22a6074e814c97c4c3d9831c3",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9306,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\italian trambling hardcore hot (!) hotel  (Curtney,Jenna).mpg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\cum [bangbus] circumcision  (Tatjana).rar.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\malaysia fetish big balls .mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\french beast beastiality licking glans .mpeg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\gang bang girls titts 40+ .avi.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\kicking licking blondie .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\tyrkish action hot (!) .mpg.exe",
        "C:\\Program Files\\dotnet\\shared\\asian kicking girls boobs ejaculation .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "beast catfight titts traffic .rar.exe",
        "beast hidden hairy  (Sonja,Curtney).avi.exe",
        "beast masturbation .mpeg.exe",
        "beast public titts .mpeg.exe",
        "black fetish fucking [milf]  (Jade).avi.exe",
        "black fetish lingerie [bangbus] titts circumcision .mpg.exe",
        "black porn sperm [free] .mpg.exe",
        "blowjob hot (!) black hairunshaved .mpeg.exe",
        "brasilian cumshot bukkake [bangbus] titts gorgeoushorny  (Jade).mpeg.exe",
        "brasilian gang bang hardcore [bangbus] cock .zip.exe",
        "danish cumshot hardcore uncut  (Sylvia).zip.exe",
        "danish handjob blowjob masturbation stockings .mpeg.exe",
        "danish kicking lesbian big pregnant .mpeg.exe",
        "danish nude lingerie [bangbus] leather .mpg.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking lesbian redhair .mpg.exe",
        "gay several models young .avi.exe",
        "horse [free] feet mistress  (Melissa).rar.exe",
        "indian action xxx [milf] circumcision .zip.exe",
        "italian cum trambling hidden leather  (Sandy,Sylvia).mpeg.exe",
        "kicking xxx hot (!) hotel .mpeg.exe",
        "lesbian hot (!) gorgeoushorny  (Ashley,Jade).avi.exe",
        "mssrv.exe",
        "rundll64.exe",
        "russian nude xxx masturbation ejaculation .rar.exe",
        "screensaver.scr",
        "setup.exe",
        "soccer.exe",
        "swedish animal lingerie masturbation .avi.exe",
        "swedish cum gay uncut titts .mpeg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}