# dataset/preprocessed/32cd02e60cfdf5932a8e8d504c352422fbc99f2d3d3689c7191ba217257e0fab.json
{
  "hash": "32cd02e60cfdf5932a8e8d504c352422fbc99f2d3d3689c7191ba217257e0fab",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  32cd02e60cfdf593_documents and settings .exe\n    C:\\ProgramData\\lvbtgx.exe\n\nCOMMAND EXECUTIONS (2):\n  (open) C:\\32cd02e60cfdf5932a8e8d504c352422fbc99f2d3d3689c7191ba217257(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\32cd02e60cfdf593_documents and settings\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - contain obfuscated stackstrings\n  - get startup folder\n  - set file attributes\n  - get common file path\n  - accept command line arguments\n\nCAPA RULE MATCHES (MBC):\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - get startup folder\n  - Process::Create Process [C0017]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - contain loop\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Delete File [C0047]\n  - File System::Create Directory [C0046]\n  - Process::Allocate Thread Local Storage [C0040]\n  - get thread local storage value\n  - Process::Set Thread Local Storage Value [C0041]\n  - File System::Read File [C0051]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - resolve function by parsing PE exports\n  - contain a resource (.rsrc) section\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  AUTOEXEC.BAT .exe  [sha256: 16676ad88b7758f6...]\n  boot.ini .exe  [sha256: c2947aa36b924faf...]\n  CONFIG.SYS .exe  [sha256: eb7497198b53ffea...]\n  cuckoo .exe  [sha256: c7ce091f9f18b492...]\n  Documents and Settings .exe  [sha256: 32a7d4cf22db03a8...]\n  eula.1028.txt .exe  [sha256: 45a94734127228ce...]\n  eula.1031.txt .exe  [sha256: 58633bb1ba0b3edb...]\n  eula.1033.txt .exe  [sha256: dbe0267c8738bbfc...]\n  eula.1036.txt .exe  [sha256: 6d0648a4910fda2c...]\n  eula.1040.txt .exe  [sha256: c4063bff2299d0b1...]\n  eula.1041.txt .exe  [sha256: ea56f24d003c15c8...]\n  eula.1042.txt .exe  [sha256: 563e2502430331d6...]\n  eula.2052.txt .exe  [sha256: 03557ba8f256bc40...]\n  eula.3082.txt .exe  [sha256: 44d24531ee8f699d...]\n  fvbro.exe  [sha256: 79947a00ada5f25c...]\n  globdata.ini .exe  [sha256: dc2f9089211e21be...]\n  install.exe .exe  [sha256: 4878b750489153df...]\n  install.ini .exe  [sha256: 6dc981290265d238...]\n  install.res.1028.dll .exe  [sha256: c6cb7898a7c3ffc9...]\n  install.res.1031.dll .exe  [sha256: d5c17ef1273cac8a...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\All Users\\Application Data\\fvbro.exe\n  C:\\ProgramData\\Saaaalamm\\Mira.h\n  c:\\$Recycle.Bin .exe\n  c:\\AUTOEXEC.BAT .exe\n  c:\\Documents and Settings .exe\n  c:\\Miras\n  c:\\PerfLogs .exe\n  c:\\Program Files (x86) .exe\n  c:\\Program Files .exe\n  c:\\ProgramData .exe\n  c:\\Recovery .exe\n  c:\\System Volume Information .exe\n  c:\\Users .exe\n  c:\\Windows .exe\n  c:\\dev_man.exe .exe\n  c:\\l.bat .exe\n  c:\\pagefile.sys .exe\n  c:\\r.py .exe\n  c:\\r.pyw .exe\n  c:\\sdel.exe .exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft\u00ae Windows\u00ae Operating System = C:\\ProgramData\\lvbtgx.exe\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System = C:\\Documents and Settings\\All Users\\Application Data\\ujedq.exe\u0000\n\nMUTEXES CREATED (7):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\All Users\\Application Data\\fvbro.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA\n\nCOMPILER AND PACKER:\n  [Compiler] MinGW\n  [Linker] GNU linker ld (GNU Binutils)\n  Packers: {'PEiD': 'MingWin32 v?.? (h)'}\n\nPE SECTIONS (entropy):\n  .text        entropy=5.66\n  .data        entropy=0.71\n  .rdata       entropy=5.01\n  .bss         entropy=0.0\n  .idata       entropy=4.29\n  .rsrc        entropy=6.06\n  [overlay]    entropy=6.987114429473877  size=417983\n",
  "ground_truth": {
    "hash": "32cd02e60cfdf5932a8e8d504c352422fbc99f2d3d3689c7191ba217257e0fab",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9167,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "get startup folder"
      }
    ],
    "capa_sig_names": [
      "contain obfuscated stackstrings",
      "get startup folder",
      "create process on Windows",
      "set file attributes",
      "delay execution",
      "create or open file",
      "contain loop",
      "get common file path",
      "delete file",
      "create directory",
      "allocate thread local storage",
      "get thread local storage value",
      "set thread local storage value",
      "read file on Windows",
      "accept command line arguments",
      "resolve function by parsing PE exports",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft\u00ae Windows\u00ae Operating System",
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "AUTOEXEC.BAT .exe",
        "C:\\ProgramData\\Saaaalamm\\Mira.h",
        "C:\\ProgramData\\lvbtgx.exe",
        "CONFIG.SYS .exe",
        "Documents and Settings .exe",
        "IO.SYS .exe",
        "Mira.h",
        "boot.ini .exe",
        "c:\\$Recycle.Bin .exe",
        "c:\\Documents and Settings .exe",
        "c:\\PerfLogs .exe",
        "c:\\Program Files (x86) .exe",
        "c:\\Program Files .exe",
        "c:\\ProgramData .exe",
        "c:\\Recovery .exe",
        "c:\\System Volume Information .exe",
        "c:\\Users .exe",
        "c:\\Windows .exe",
        "c:\\dev_man.exe .exe",
        "c:\\l.bat .exe",
        "c:\\pagefile.sys .exe",
        "c:\\r.py .exe",
        "c:\\r.pyw .exe",
        "c:\\sdel.exe .exe",
        "cuckoo .exe",
        "eula.1028.txt .exe",
        "eula.1031.txt .exe",
        "eula.1033.txt .exe",
        "eula.1036.txt .exe",
        "eula.1040.txt .exe",
        "eula.1041.txt .exe",
        "eula.1042.txt .exe",
        "eula.2052.txt .exe",
        "eula.3082.txt .exe",
        "fvbro.exe",
        "globdata.ini .exe",
        "install.exe .exe",
        "install.ini .exe",
        "install.res.1028.dll .exe",
        "install.res.1031.dll .exe",
        "install.res.1033.dll .exe",
        "install.res.1036.dll .exe",
        "install.res.1040.dll .exe",
        "install.res.1041.dll .exe",
        "install.res.1042.dll .exe",
        "install.res.2052.dll .exe",
        "install.res.3082.dll .exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\Documents and Settings\\All Users\\Application Data\\fvbro.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}