# dataset/preprocessed/3321ab16cc95be75538a0318ada5cf61336da24231dfa4bad754921fed119cb3.json
{
  "hash": "3321ab16cc95be75538a0318ada5cf61336da24231dfa4bad754921fed119cb3",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (9):\n  GET http://carder.bit/  UA: Mozilla/5.0 (Windows NT 6.1; WOW64; Trident/7.0; rv:11.0) like Gecko\n  GET http://carder.bit/\n  GET http://ipv4bot.whatismyipaddress.com/\n  GET http://110.110.110.0/wpad.dat\n  POST http://110.110.110.3/aubey?zaow=pl&lf=orestore\n  POST http://110.110.110.3/plauza?pley=ghea&eage=ee\n  POST http://110.110.110.3/zaeeza?ee=orephai&eapl=stow\n  POST http://110.110.110.3/ieplui\n  POST http://110.110.110.3/zaea?scast=eiplai\n\nDNS LOOKUPS (8):\n  carder.bit\n  ipv4bot.whatismyipaddress.com\n  ipv4bot.whatismyipaddress.com \u2192 66.171.248.178\n  ns1.wowservers.ru\n  ns1.wowservers.ru \u2192 185.103.243.59\n  ns2.wowservers.ru\n  ns2.wowservers.ru \u2192 185.103.243.59\n  ransomware.bit\n\nIP TRAFFIC (2):\n  185.103.243.59:53 (UDP)\n  66.171.248.178:80 (TCP)\n\nSERVICES OPENED (2):\n  RASMAN\n  Sens\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  C:\\Users\\Lucas\\AppData\\Local\\Temp\\EY0d6zHPqPMtxB.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n\nCOMMAND EXECUTIONS (5):\n  C:\\Users\\Lucas\\AppData\\Local\\Temp\\EY0d6zHPqPMtxB.exe\n  nslookup carder.bit ns1.wowservers.ru\n  nslookup ransomware.bit ns2.wowservers.ru\n  nslookup carder.bit ns2.wowservers.ru\n  nslookup ransomware.bit ns1.wowservers.ru\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (6):\n  - contain obfuscated stackstrings\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - modify service\n  - open clipboard\n  - get common file path\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Memory::Allocate Memory [C0007]\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - contain loop\n  - Process::Terminate Process [C0018]\n  - Process::Allocate Thread Local Storage [C0040]\n  - link function at runtime on Windows\n  - link many functions at runtime\n  - modify service\n  - open clipboard\n  - Discovery::File and Directory Discovery [E1083]\n  - linked against CPP standard library\n  - contain a resource (.rsrc) section\n\nHIGHLIGHTED CALLS:\n  SetFileTime\n\n--- ARTIFACTS ---\n\nFILES DROPPED (2):\n  seoljp.exe  [sha256: d53b7fb20e7cc440...]\n  C:\\Users\\<USER>\\AppData\\Roaming\\Microsoft\\gghdmp.exe  [sha256: bf13e9164fd99733...]\n\nFILES WRITTEN (3):\n  C:\\Documents and Settings\\<USER>\\Application Data\\Microsoft\\seoljp.exe\n  C:\\Documents and Settings\\Administrator\\Application Data\\Microsoft\\mnesdv.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Microsoft\\Windows\\Temporary Internet Files\\Content.IE5\\XBL73YY1\\ipv4bot_whatismyipaddress_com[1].htm\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\rihpwflfxqa = \"\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\nonsjazdbqw = \"C:\\Documents and Settings\\Administrator\\Application Data\\Microsoft\\mnesdv.exe\"\u0000\n\nOTHER REGISTRY KEYS SET (10):\n  HKU\\S-1-5-18\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\ProxyEnable\n  HKU\\S-1-5-18\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\SavedLegacySettings\n  HKU\\S-1-5-18\\Software\\Microsoft\\windows\\CurrentVersion\\Internet Settings\\Connections\\DefaultConnectionSettings\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\ZONEMAP\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\WPAD\\52-55-\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\5.0\\CACHE\\H\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\CONNECTIONS\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\WPAD\\{3903E\n  HKU\\S-1-5-21-3712457824-2419000099-45725732-1005\\SOFTWARE\\MICROSOFT\\WINDOWS\\CURRENTVERSION\\INTERNET SETTINGS\\5.0\\CACHE\\C\n\nMUTEXES CREATED (17):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  Global\\pc_group=WORKGROUP&ransom_id=2d84a7bff44048e7\n  Global\\pc_group=WORKGROUP&ransom_id=623e60a7ec07c162\n  Global\\pc_group=WORKGROUP&ransom_id=ea2177a071d2c847\n  IESQMMUTEX_0_208\n\nMUTEXES OPENED (9):\n  Global\\pc_group=WORKGROUP&ransom_id=cd19fdb96c78a9c3\n  Local\\MSCTF.Asm.MutexDefault1\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!history!history.ie5!\n  Local\\c:!users!<USER>!appdata!local!microsoft!windows!temporary internet files!content.ie5!\n  Local\\c:!users!<USER>!appdata!roaming!microsoft!windows!cookies!\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\nPROCESSES CREATED (4):\n  C:\\Users\\Lucas\\AppData\\Local\\Temp\\EY0d6zHPqPMtxB.exe\n  C:\\WINDOWS\\system32\\nslookup.exe\n  C:\\Windows\\SysWOW64\\nslookup.exe\n  nslookup carder.bit ns1.wowservers.ru\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, LoadLibraryW, VirtualProtect\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n\nPE SECTIONS (entropy):\n  .text        entropy=6.61\n  .rdata       entropy=5.83\n  .data        entropy=4.12\n  .rsrc        entropy=7.8\n  .reloc       entropy=1.3\n  [overlay]    entropy=2.503258466720581  size=9\n",
  "ground_truth": {
    "hash": "3321ab16cc95be75538a0318ada5cf61336da24231dfa4bad754921fed119cb3",
    "verdict": "MALICIOUS",
    "detection_rate": 0.942,
    "malicious_count": 65,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "ransom",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1115",
        "signature_description": "open clipboard"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1543",
        "signature_description": "modify service"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "allocate memory",
      "contain obfuscated stackstrings",
      "contain loop",
      "terminate process",
      "allocate thread local storage",
      "link function at runtime on Windows",
      "link many functions at runtime",
      "modify service",
      "open clipboard",
      "get common file path",
      "linked against CPP standard library",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [
        "185.103.243.59",
        "66.171.248.178"
      ],
      "domains": [
        "carder.bit",
        "ipv4bot.whatismyipaddress.com",
        "ns1.wowservers.ru",
        "ns2.wowservers.ru",
        "ransomware.bit"
      ],
      "urls": [
        "http://110.110.110.0/wpad.dat",
        "http://110.110.110.3/aubey?zaow=pl&lf=orestore",
        "http://110.110.110.3/ieplui",
        "http://110.110.110.3/plauza?pley=ghea&eage=ee",
        "http://110.110.110.3/zaea?scast=eiplai",
        "http://110.110.110.3/zaeeza?ee=orephai&eapl=stow",
        "http://carder.bit/",
        "http://ipv4bot.whatismyipaddress.com/"
      ],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\rihpwflfxqa",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\nonsjazdbqw"
      ],
      "registry_all_count": 16,
      "files_dropped": [
        "C:\\Users\\<USER>\\AppData\\Roaming\\Microsoft\\gghdmp.exe",
        "seoljp.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "Global\\pc_group=WORKGROUP&ransom_id=2d84a7bff44048e7",
        "Global\\pc_group=WORKGROUP&ransom_id=623e60a7ec07c162",
        "Global\\pc_group=WORKGROUP&ransom_id=ea2177a071d2c847",
        "IESQMMUTEX_0_208",
        "Local\\ZonesCacheCounterMutex",
        "Local\\ZonesLockedCacheCounterMutex",
        "RasPbFile",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "C:\\Users\\Lucas\\AppData\\Local\\Temp\\EY0d6zHPqPMtxB.exe",
        "C:\\WINDOWS\\system32\\nslookup.exe",
        "C:\\Windows\\SysWOW64\\nslookup.exe",
        "nslookup carder.bit ns1.wowservers.ru"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Lastline",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}