# dataset/preprocessed/33673f0d3bef664821bad28ab32a7a497be46e796f8eabb91460172d45663216.json
{
  "hash": "33673f0d3bef664821bad28ab32a7a497be46e796f8eabb91460172d45663216",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  bukkake-girls-mpeg.pe32\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n    rundll32.exe\n\nCOMMAND EXECUTIONS (3):\n  ((null)) C:\\33673f0d3bef664821bad28ab32a7a497be46e796f8eabb91460172d45663216(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\bukkake girls .mpeg.exe\n  \"C:\\Windows\\system32\\rundll32.exe\" C:\\Windows\\system32\\shell32.dll,OpenAs_RunDLL C:\\Users\\Petra\\AppData\\Local\\Temp\\bukkake-girls-mpeg.pe32\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american gang bang xxx sleeping blondie .avi.exe  [sha256: 0112cc0484d7ba6b...]\n  beast full movie  (Tatjana).rar.exe  [sha256: b00abfb153c9edf8...]\n  black cumshot bukkake hidden wifey .rar.exe  [sha256: 7190c1944ea8b967...]\n  black fetish hardcore uncut .zip.exe  [sha256: bc777228ddbbac74...]\n  black gang bang beast lesbian feet bondage .mpg.exe  [sha256: a42822a6e95f7da6...]\n  black kicking bukkake uncut hole .mpeg.exe  [sha256: 4bbccc2ebb8a08c2...]\n  blowjob hidden redhair  (Jenna,Janette).rar.exe  [sha256: 2f8aecc4548e1409...]\n  blowjob public feet .rar.exe  [sha256: cc28002bbbac9cb6...]\n  brasilian porn trambling sleeping  (Karin).zip.exe  [sha256: cdbb7988000aa94a...]\n  danish cum trambling full movie sweet .avi.exe  [sha256: 7aa288abd760b94c...]\n  danish handjob xxx hot (!) cock  (Sonja,Samantha).mpg.exe  [sha256: 62e01bd30fc89db5...]\n  C:\\Windows\\mssrv.exe  [sha256: 84622caca9993a2f...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\black fetish lesbian [free] mature .avi.exe  [sha256: 52ccaa61439d42b6...]\n  C:\\Program Files\\dotnet\\shared\\indian animal sperm [bangbus] black hairunshaved .avi.exe  [sha256: f3115474c361c259...]\n  C:\\Program Files\\DVD Maker\\Shared\\brasilian action hardcore voyeur  (Liz).avi.exe  [sha256: f5d62c28339755bf...]\n  C:\\Program Files\\Windows Journal\\Templates\\italian nude trambling [milf] cock balls  (Melissa).zip.exe  [sha256: f2ccc71eb8e72011...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\italian handjob horse catfight hairy .rar.exe  [sha256: 408b3920d58afbfb...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\sperm [free] traffic .rar.exe  [sha256: 965af5acc588e4c2...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\russian horse sperm several models balls  (Gina,Sylvia).mpeg.exe  [sha256: e1b7cd674f659bd4...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\american cumshot trambling masturbation .avi.exe  [sha256: bfcda136f3a05c37...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\black cumshot bukkake hidden wifey .rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\danish cum trambling full movie sweet .avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\indian horse bukkake licking hole .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\russian gang bang xxx hidden  (Karin).avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\danish handjob xxx hot (!) cock  (Sonja,Samantha).mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\italian nude blowjob catfight glans mature .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\hardcore public glans pregnant .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\american gang bang xxx sleeping blondie .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\fucking [milf]  (Tatjana).mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\beast full movie  (Tatjana).rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\japanese cumshot lesbian hot (!) pregnant  (Sonja,Karin).rar.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\blowjob hidden redhair  (Jenna,Janette).rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\lesbian hidden hole .mpeg.exe\n  C:\\WINDOWS\\Temp\\brasilian porn trambling sleeping  (Karin).zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\lingerie several models .rar.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\blowjob public feet .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\russian action xxx hot (!) stockings .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\swedish handjob blowjob hot (!) hole mistress  (Melissa).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\black kicking bukkake uncut hole .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nOTHER REGISTRY KEYS SET (1):\n  \\REGISTRY\\MACHINE\\HexStream = 53004F004600540057004100520045005C004D006900630072006F0073006F00660074005C00570069006E0064\n\nMUTEXES CREATED (1):\n  mutex666\n\nPROCESSES CREATED (1):\n  rundll32.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880704402923584  size=129785\n",
  "ground_truth": {
    "hash": "33673f0d3bef664821bad28ab32a7a497be46e796f8eabb91460172d45663216",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\sperm [free] traffic .rar.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\russian horse sperm several models balls  (Gina,Sylvia).mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\american cumshot trambling masturbation .avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\sperm [milf] .zip.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\japanese cumshot lingerie [milf] hole bedroom  (Karin).avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\beast sleeping hole sm .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\american nude lingerie [bangbus] feet .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\russian kicking lesbian catfight glans castration .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\danish gang bang trambling uncut hole traffic  (Curtney).avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\black fetish lesbian [free] mature .avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\brasilian action hardcore voyeur  (Liz).avi.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\italian nude trambling [milf] cock balls  (Melissa).zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\italian handjob horse catfight hairy .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\indian animal sperm [bangbus] black hairunshaved .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "\\Program Files\\Common Files\\microsoft shared\\action lesbian granny .zip.exe",
        "\\Program Files\\DVD Maker\\Shared\\brasilian beastiality horse voyeur boots .mpeg.exe",
        "\\Program Files\\Google\\Update\\1.3.34.11\\Recovery\\GUR4F00.tmp\\french lingerie masturbation shower .zip.exe",
        "\\Program Files\\Google\\Update\\Download\\american cumshot kicking hot (!) boobs femdom  (Karin,Curtney).mpg.exe",
        "\\Program Files\\Microsoft Office\\Templates\\chinese animal horse sleeping  (Sonja).zip.exe",
        "\\Program Files\\Windows Journal\\Templates\\lingerie licking vagina femdom .mpeg.exe",
        "\\Program Files\\Windows Sidebar\\Shared Gadgets\\kicking lesbian titts redhair .mpeg.exe",
        "\\ProgramData\\Microsoft\\Network\\Downloader\\canadian porn beastiality voyeur beautyfull .avi.exe",
        "\\Windows\\mssrv.exe",
        "american gang bang xxx sleeping blondie .avi.exe",
        "beast full movie  (Tatjana).rar.exe",
        "black cumshot bukkake hidden wifey .rar.exe",
        "black fetish hardcore uncut .zip.exe",
        "black gang bang beast lesbian feet bondage .mpg.exe",
        "black kicking bukkake uncut hole .mpeg.exe",
        "blowjob hidden redhair  (Jenna,Janette).rar.exe",
        "blowjob public feet .rar.exe",
        "brasilian porn trambling sleeping  (Karin).zip.exe",
        "danish cum trambling full movie sweet .avi.exe",
        "danish handjob xxx hot (!) cock  (Sonja,Samantha).mpg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": [
        "rundll32.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "SNDBOX",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}