# dataset/preprocessed/339d5941291b277feb92ce5d8796611059c50579aa9f936723dc08320c3fc52b.json
{
  "hash": "339d5941291b277feb92ce5d8796611059c50579aa9f936723dc08320c3fc52b",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  339d5941291b277f_xxx licking glans  (anniston,janette).mpg.exe\n    C:\\Users\\<USER>\\Downloads\\339d5941291b277f_xxx licking glans  (anniston,janette).mpg.exe\n    C:\\Users\\<USER>\\Downloads\\339d5941291b277f_xxx licking glans  (anniston,janette).mpg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\339d5941291b277feb92ce5d8796611059c50579aa9f936723dc08320c3fc52b(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\339d5941291b277f_xxx licking glans  (anniston,janette).mpg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american porn blowjob several models  (Karin).mpeg.exe  [sha256: 670be86fe4456201...]\n  blowjob big penetration .mpeg.exe  [sha256: e6e4d0af3af6bb46...]\n  brasilian porn bukkake [milf]  (Sylvia).rar.exe  [sha256: 8cef5907d2805095...]\n  bukkake sleeping fishy .mpeg.exe  [sha256: 7939a6cb5680fff2...]\n  danish horse horse [bangbus] .mpeg.exe  [sha256: 2ec0e470975d7eea...]\n  danish porn hardcore big glans .mpeg.exe  [sha256: 0263ee4d860013bb...]\n  defrag.exe  [sha256: 79b84c847a479039...]\n  defrag64.exe  [sha256: 462a01d5843204e6...]\n  driver.exe  [sha256: 2b2d19486e421075...]\n  drwatson32.exe  [sha256: c078d31bd205ec7e...]\n  fucking sleeping cock .zip.exe  [sha256: 06776f15b3ec1006...]\n  gay catfight glans .avi.exe  [sha256: 1264f52282500862...]\n  IE9setup.exe  [sha256: 91b3df255713b7a9...]\n  indian action gay [milf] titts .zip.exe  [sha256: 65df1cc16eaefbb5...]\n  indian action trambling licking hole traffic  (Melissa).mpg.exe  [sha256: 63ffc7aa95ab38b4...]\n  indian handjob horse big glans .mpg.exe  [sha256: 18acdd45f3f28cd4...]\n  indian handjob horse masturbation .avi.exe  [sha256: 19c67e8c5f6833b7...]\n  indian kicking lesbian licking granny .zip.exe  [sha256: 62157e11299283c7...]\n  indian kicking xxx licking cock latex .mpeg.exe  [sha256: a7f86ad4340efe81...]\n  italian beastiality sperm lesbian feet traffic  (Curtney).avi.exe  [sha256: bd1e07d4926e9001...]\n\nFILES WRITTEN (30):\n  C:\\339d5941291b277feb92ce5d8796611059c50579aa9f936723dc08320c3fc52b\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\japanese kicking gay big hole .mpg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\indian handjob horse masturbation .avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\indian action gay [milf] titts .zip.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\japanese horse trambling public .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\danish horse horse [bangbus] .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\spanish hardcore hot (!) glans girly  (Karin).rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\indian kicking lesbian licking granny .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\sperm several models 40+ .zip.exe\n  C:\\Program Files\\Movie Maker\\Shared\\russian beastiality lesbian catfight ejaculation  (Sonja,Liz).avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\indian handjob horse big glans .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\bukkake sleeping fishy .mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\gay catfight glans .avi.exe\n  C:\\WINDOWS\\Temp\\sperm several models titts high heels .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\sperm public glans  (Christine,Sarah).mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\danish porn hardcore big glans .mpeg.exe\n  C:\\WINDOWS\\assembly\\temp\\swedish fetish gay licking balls .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\indian action trambling licking hole traffic  (Melissa).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\italian beastiality sperm lesbian feet traffic  (Curtney).avi.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.879477500915527  size=920561\n",
  "ground_truth": {
    "hash": "339d5941291b277feb92ce5d8796611059c50579aa9f936723dc08320c3fc52b",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\american hardcore sleeping .avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\brasilian gang bang [bangbus] sm .avi.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\norwegian lingerie hidden .mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\gay uncut .zip.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\lingerie gang bang catfight .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\swedish blowjob hot (!) femdom .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\chinese handjob hardcore [bangbus] legs traffic  (Jenna).zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\gang bang voyeur titts redhair .mpg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\handjob hardcore several models glans gorgeoushorny  (Janette,Sylvia).rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\canadian fetish sleeping boobs .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\german cumshot catfight .avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\danish beast voyeur .rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\american cum lesbian big  (Tatjana,Jade).mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\porn hidden bedroom .zip.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\fucking masturbation legs .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\trambling handjob big swallow .zip.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american porn blowjob several models  (Karin).mpeg.exe",
        "blowjob big penetration .mpeg.exe",
        "brasilian porn bukkake [milf]  (Sylvia).rar.exe",
        "bukkake sleeping fishy .mpeg.exe",
        "danish horse horse [bangbus] .mpeg.exe",
        "danish porn hardcore big glans .mpeg.exe",
        "defrag.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking sleeping cock .zip.exe",
        "gay catfight glans .avi.exe",
        "indian action gay [milf] titts .zip.exe",
        "indian action trambling licking hole traffic  (Melissa).mpg.exe",
        "indian handjob horse big glans .mpg.exe",
        "indian handjob horse masturbation .avi.exe",
        "indian kicking lesbian licking granny .zip.exe",
        "indian kicking xxx licking cock latex .mpeg.exe",
        "italian beastiality sperm lesbian feet traffic  (Curtney).avi.exe",
        "italian cumshot beast [milf]  (Janette).zip.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}