# dataset/preprocessed/3523161b56d9b37e7f634860ce2f52df478e249f8be17eb830cfd6a017308cbe.json
{
  "hash": "3523161b56d9b37e7f634860ce2f52df478e249f8be17eb830cfd6a017308cbe",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  ftp.byethost12.com\n\nSERVICES OPENED (2):\n  RASMAN\n  eventchk\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  services.exe\n    svchost.exe\n  ****.exe\n    jusched.exe\n\nCOMMAND EXECUTIONS (1):\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (21):\n  - set file attributes\n  - parse PE header\n  - get file size\n  - get common file path\n  - write process memory\n  - link function at runtime on Windows\n  - check if file exists\n  - enumerate files on Windows\n  - enumerate files recursively\n  - get disk information\n  - query or enumerate registry value\n  - get socket status\n  - download and write a file\n  - query service status\n  - create service\n  - persist via Windows service\n  - start service\n  - modify service\n  - delete registry key\n  - enumerate processes\n  - link many functions at runtime\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - File System::Set File Attributes [C0050]\n  - File System::Create File [C0016]\n  - Process::Create Process [C0017]\n  - Process::Create Thread [C0038]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Get File Attributes [C0049]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - contain loop\n  - parse PE header\n  - resolve function by parsing PE exports\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - Execution::Install Additional Program [B0023]\n  - get OS version\n  - Communication::Socket Communication::Initialize Winsock Library [C0001.009]\n  - write process memory\n  - link function at runtime on Windows\n  - Discovery::Code Discovery::Enumerate PE Sections [B0046.001]\n  - File System::Create Directory [C0046]\n  - File System::Delete File [C0047]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (1):\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: 086d9176bd870192...]\n\nFILES WRITTEN (8):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files\\f440491e\\f440491e\n  C:\\Program Files\\f440491e\\jusched.exe\n  C:\\WINDOWS\\Tasks\\Update23.job\n  C:\\WINDOWS\\system32\\inf\\svchost.exe\n  C:\\Windows\\INETINFO.exe\n  C:\\Windows\\Tasks\\Update23.job\n  C:\\Windows\\system32\\inf\\svchost.exe\n\nPERSISTENCE REGISTRY KEYS (4):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:JavaUpdate23\n  HKLM\\system\\controlset001\\services\\sharedaccess\\parameters\\firewallpolicy\\standardprofile\\authorizedapplications\\list\\C:\\Windows\\system32\\inf\\svchost.exe = C:\\Windows\\system32\\inf\\svchost.exe:*:Enabled:@xpsp2res.dll,-22001\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet001\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\WINDOWS\\system32\\inf\\svchost.exe = C:\\WINDOWS\\system32\\inf\\svchost.exe:*:Enabled:@xpsp2res.dll,-22001\u0000\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe = C:\\Program Files\\f440491e\\jusched.exe:*:Enabled:JavaUpdate23\u0000\n\nOTHER REGISTRY KEYS SET (2):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Program \n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  RasPbFile\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n  ZonesLockedCacheCounterMutex\n\nMUTEXES OPENED (4):\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc, WriteProcessMemory\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n\nPE SECTIONS (entropy):\n  .text        entropy=3.88\n  .rdata       entropy=4.35\n  .data        entropy=0.44\n  .idata       entropy=2.9\n  .rsrc        entropy=2.5\n  .reloc       entropy=4.41\n  [overlay]    entropy=7.076350212097168  size=11723\n",
  "ground_truth": {
    "hash": "3523161b56d9b37e7f634860ce2f52df478e249f8be17eb830cfd6a017308cbe",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9296,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "query service status"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1016",
        "signature_description": "get socket status"
      },
      {
        "id": "T1055",
        "signature_description": "write process memory"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get file size"
      },
      {
        "id": "T1105",
        "signature_description": "download and write a file"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry key"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "create service"
      },
      {
        "id": "T1569",
        "signature_description": "create service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "set file attributes",
      "create or open file",
      "create process on Windows",
      "create thread",
      "delay execution",
      "get file attributes",
      "create or open registry key",
      "create UDP socket",
      "contain loop",
      "parse PE header",
      "resolve function by parsing PE exports",
      "get file size",
      "read file via mapping",
      "write file on Windows",
      "write and execute a file",
      "get common file path",
      "get OS version",
      "initialize Winsock library",
      "write process memory",
      "link function at runtime on Windows",
      "enumerate PE sections",
      "create directory",
      "delete file",
      "read file on Windows",
      "check if file exists",
      "set registry value",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get disk information",
      "query or enumerate registry value",
      "send data on socket",
      "send data",
      "read and send data from client to server",
      "resolve DNS",
      "get socket status",
      "set socket configuration",
      "receive data on socket",
      "receive data",
      "download and write a file",
      "receive and write data from server to client",
      "get service handle",
      "query service status",
      "create service",
      "persist via Windows service",
      "start service",
      "modify service",
      "delete registry key",
      "enumerate processes",
      "terminate process",
      "allocate thread local storage",
      "link many functions at runtime",
      "contain an embedded PE file",
      "contain a resource (.rsrc) section",
      "debug build",
      "run as service"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "ftp.byethost12.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "HKLM\\system\\controlset001\\services\\sharedaccess\\parameters\\firewallpolicy\\standardprofile\\authorizedapplications\\list\\C:\\Windows\\system32\\inf\\svchost.exe",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet001\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\WINDOWS\\system32\\inf\\svchost.exe",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe"
      ],
      "registry_all_count": 6,
      "files_dropped": [
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "RasPbFile",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Jujubox"
  ]
}