# dataset/preprocessed/3612db525e87dd5808ee96772c8e069ea12aff6993021d06fa5e1cf2a0704ae9.json
{
  "hash": "3612db525e87dd5808ee96772c8e069ea12aff6993021d06fa5e1cf2a0704ae9",
  "report_text": "--- NETWORK ACTIVITY ---\n\nHTTP CONVERSATIONS (7):\n  GET http://fe.trymedia.com/a.a  UA: Macrovision_DM_2.2.146\n  GET http://www.fe.trymedia.com/a.a  UA: Macrovision_DM_2.2.146\n  GET http://www.w3.org/TR/xhtml1/DTD/xhtml1-transitional.dtd  UA: Mozilla/4.0 (compatible; MSIE 6.0; Windows NT 5.1; SV1; .NET CLR 2.0.50727; .NET CLR 3.0.04506.648; .NET CLR 3.5.21022)\n  GET http://www.w3.org/TR/xhtml1/DTD/xhtml-lat1.ent  UA: Mozilla/4.0 (compatible; MSIE 6.0; Windows NT 5.1; SV1; .NET CLR 2.0.50727; .NET CLR 3.0.04506.648; .NET CLR 3.5.21022)\n  GET http://www.w3.org/TR/xhtml1/DTD/xhtml-symbol.ent  UA: Mozilla/4.0 (compatible; MSIE 6.0; Windows NT 5.1; SV1; .NET CLR 2.0.50727; .NET CLR 3.0.04506.648; .NET CLR 3.5.21022)\n  GET http://www.w3.org/TR/xhtml1/DTD/xhtml-special.ent  UA: Mozilla/4.0 (compatible; MSIE 6.0; Windows NT 5.1; SV1; .NET CLR 2.0.50727; .NET CLR 3.0.04506.648; .NET CLR 3.5.21022)\n  GET http://fe.trymedia.com/a.a\n\nDNS LOOKUPS (4):\n  fe.trymedia.com\n  fe.trymedia.com \u2192 149.13.34.195\n  www.fe.trymedia.com \u2192 149.13.34.195\n  www.w3.org \u2192 128.30.52.100\n\nIP TRAFFIC (3):\n  128.30.52.100:80 (TCP)\n  149.13.34.195:80 (TCP)\n  <MACHINE_DNS_SERVER>:53 (UDP)\n\nSERVICES OPENED (1):\n  RASMAN\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (17):\n  - check Internet connectivity via WinINet\n  - hide graphical window\n  - parse PE header\n  - check OS version\n  - query or enumerate registry value\n  - link function at runtime on Windows\n  - get disk size\n  - find graphical window\n  - get common file path\n  - persist via Run registry key\n  - delete registry value\n  - delete registry key\n  - get file size\n  - download and write a file\n  - accept command line arguments\n  - query or enumerate registry key\n  - link many functions at runtime\n\nCAPA RULE MATCHES (MBC):\n  - File System::Create File [C0016]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - File System::Get File Attributes [C0049]\n  - parse URL\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - reference absolute stream path on Windows\n  - Communication::HTTP Communication [C0002]\n  - Process::Create Process [C0017]\n  - Process::Create Thread [C0038]\n  - Process::Terminate Thread [C0039]\n  - check Internet connectivity via WinINet\n  - hide graphical window\n  - Memory::Allocate Memory [C0007]\n  - File System::Read File [C0051]\n  - parse PE header\n  - contain loop\n  - File System::Writes File [C0052]\n  - get OS version\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - link function at runtime on Windows\n  - print debug messages\n  - Process::Create Mutex [C0042]\n  - find graphical window\n\n--- ARTIFACTS ---\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\00000000000000000000000000000000 = C:\\DOCUME~1\\ADMINI~1\\LOCALS~1\\Temp\\EB93A6\\996E.exe /r\u0000\n\nOTHER REGISTRY KEYS SET (5):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Internet Explorer\\International\\CpMRU\\En\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Internet Explorer\\International\\CpMRU\\Si\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Internet Explorer\\International\\CpMRU\\In\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Internet Explorer\\International\\CpMRU\\Fa\n\nMUTEXES CREATED (11):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  MSCTF.Shared.MUTEX.EBH\n  RasPbFile\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n\nMUTEXES OPENED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1275210071-920026266-1060284298-1003\n  CTF.Compart.MutexDefaultS-1-5-21-1275210071-920026266-1060284298-1003\n  CTF.LBES.MutexDefaultS-1-5-21-1275210071-920026266-1060284298-1003\n  CTF.Layouts.MutexDefaultS-1-5-21-1275210071-920026266-1060284298-1003\n  CTF.TMD.MutexDefaultS-1-5-21-1275210071-920026266-1060284298-1003\n  MSIMGSIZECacheMutex\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, LoadLibraryW, VirtualAlloc, VirtualProtect\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Sign tool] Windows Authenticode\n\nPE SECTIONS (entropy):\n  .text        entropy=6.67\n  .rdata       entropy=4.65\n  .data        entropy=4.28\n  .rsrc        entropy=6.01\n  [overlay]    entropy=5.815910816192627  size=3912\n",
  "ground_truth": {
    "hash": "3612db525e87dd5808ee96772c8e069ea12aff6993021d06fa5e1cf2a0704ae9",
    "verdict": "MALICIOUS",
    "detection_rate": 0.6957,
    "malicious_count": 48,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "trymedia",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1016",
        "signature_description": "check Internet connectivity via WinINet"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "check OS version"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1105",
        "signature_description": "download and write a file"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry value"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      }
    ],
    "capa_sig_names": [
      "create or open file",
      "create or open registry key",
      "get file attributes",
      "parse URL",
      "delay execution",
      "reference absolute stream path on Windows",
      "get HTTP content length",
      "create process on Windows",
      "create thread",
      "terminate thread",
      "check Internet connectivity via WinINet",
      "hide graphical window",
      "allocate memory",
      "allocate RWX memory",
      "read file on Windows",
      "parse PE header",
      "contain loop",
      "write file on Windows",
      "get OS version",
      "check OS version",
      "query or enumerate registry value",
      "link function at runtime on Windows",
      "print debug messages",
      "create mutex",
      "get disk size",
      "find graphical window",
      "set registry value",
      "get common file path",
      "delete file",
      "copy file",
      "persist via Run registry key",
      "delete registry value",
      "create directory",
      "send HTTP request",
      "create HTTP request",
      "check HTTP status code",
      "connect to HTTP server",
      "send data",
      "delete registry key",
      "get file size",
      "connect to URL",
      "check for time delay via GetTickCount",
      "read data from Internet",
      "receive data",
      "download and write a file",
      "receive and write data from server to client",
      "move file",
      "resolve DNS",
      "accept command line arguments",
      "query or enumerate registry key",
      "extract resource via kernel32 functions",
      "get graphical window text",
      "check mutex and exit",
      "terminate process",
      "allocate thread local storage",
      "link many functions at runtime",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [
        "128.30.52.100",
        "149.13.34.195",
        "<MACHINE_DNS_SERVER>"
      ],
      "domains": [
        "fe.trymedia.com",
        "www.fe.trymedia.com",
        "www.w3.org"
      ],
      "urls": [
        "http://fe.trymedia.com/a.a",
        "http://www.fe.trymedia.com/a.a",
        "http://www.w3.org/TR/xhtml1/DTD/xhtml-lat1.ent",
        "http://www.w3.org/TR/xhtml1/DTD/xhtml-special.ent",
        "http://www.w3.org/TR/xhtml1/DTD/xhtml-symbol.ent",
        "http://www.w3.org/TR/xhtml1/DTD/xhtml1-transitional.dtd"
      ],
      "registry_persistence": [
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\00000000000000000000000000000000"
      ],
      "registry_all_count": 6,
      "files_dropped": [],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "MSCTF.Shared.MUTEX.EBH",
        "RasPbFile",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}