# dataset/preprocessed/3657e8733cdea568bc631f2df3ec1bae1858e31e5c88cabb84558ce36bc42094.json
{
  "hash": "3657e8733cdea568bc631f2df3ec1bae1858e31e5c88cabb84558ce36bc42094",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  <PATH_SAMPLE.EXE>\n    %WINDIR%\\syswow64\\cmd.exe\n  %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\explorer.exe\n    %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\reg.exe\n    %WINDIR%\\syswow64\\xcopy.exe\n  %WINDIR%\\syswow64\\reg.exe\n  %WINDIR%\\syswow64\\cmd.exe\n  %WINDIR%\\syswow64\\cmd.exe\n  %WINDIR%\\syswow64\\explorer.exe\n  %WINDIR%\\syswow64\\cmd.exe\n\nCOMMAND EXECUTIONS (13):\n  (open) C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\3.tmp\\promo.bat [C:\\WINDOWS\\system32\\]\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\6D09.tmp\\promo.bat\n  cmd /c \"\"C:\\Documents and Settings\\Administrator\\Local Settings\\Temp\\2.tmp\\promo.bat\" \"\n  C:\\WINDOWS\\system32\\cmd.exe /c echo 2019-11-17 \u661f\u671f\u65e5\n  C:\\WINDOWS\\system32\\cmd.exe /c echo 18:54:10.22\n  explorer  Promo\\\n  C:\\WINDOWS\\system32\\cmd.exe /c echo 18:54:10.54\n  C:\\WINDOWS\\system32\\cmd.exe /c type \"C:\\Documents and Settings\\Administrator\\Promo_time\"\n  C:\\WINDOWS\\system32\\REG ADD HKEY_LOCAL_MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run /v \"Promo\" /t REG_SZ /d \"C:\\Documents and Settings\\Administrator\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\Promo.exe\" /f\n  xcopy \"Promo\" /y /e /h \"C:\\Documents and Settings\\Administrator\\Promo\\\"\n  C:\\WINDOWS\\system32\\cmd.exe /c echo 2019-11-30 \u661f\u671f\u516d\n  C:\\WINDOWS\\system32\\cmd.exe /c echo 16:25:37.54\n  C:\\WINDOWS\\system32\\cmd.exe /c echo 16:25:37.82\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (10):\n  - encode data using XOR\n  - identify system language via API\n  - get common file path\n  - link function at runtime on Windows\n  - enumerate gui resources\n  - accept command line arguments\n  - check OS version\n  - log keystrokes via polling\n  - check for unmoving mouse cursor\n  - reference startup folder\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Process::Create Process [C0017]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - Process::Allocate Thread Local Storage [C0040]\n  - contain loop\n  - identify system language via API\n  - File System::Delete Directory [C0048]\n  - Discovery::File and Directory Discovery [E1083]\n  - extract resource via kernel32 functions\n  - link function at runtime on Windows\n  - enumerate gui resources\n  - Execution::Command and Scripting Interpreter [E1059]\n  - Process::Terminate Process [C0018]\n  - get OS version\n  - Discovery::System Information Discovery [E1082]\n  - Discovery::Application Window Discovery [E1010]\n  - Collection::Keylogging::Polling [F0002.002]\n  - Anti-Behavioral Analysis::Virtual Machine Detection::Human User Check [B0009.012]\n  - File System::Delete File [C0047]\n  - File System::Create Directory [C0046]\n  - File System::Writes File [C0052]\n  - File System::Read File [C0051]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (1):\n  %TEMP%\\61bb.tmp\\promo.bat  [sha256: 85a6032b65286a2b...]\n\nFILES WRITTEN (2):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\3.tmp\\promo.bat\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\6D09.tmp\\promo.bat\n\nPERSISTENCE REGISTRY KEYS (3):\n  <HKLM>\\SOFTWARE\\Wow6432Node\\Microsoft\\Windows\\CurrentVersion\\Run\\Promo = %APPDATA%\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\Promo.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Promo = C\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Promo = C:\\Documents and Settings\\Administrator\\AppData\\Roaming\\Microsoft\\Windows\\Start Menu\\Programs\\Startu\n\nOTHER REGISTRY KEYS SET (3):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Document\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\UNCAsIntranet\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\AutoDetect\n\nMUTEXES CREATED (11):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  MSCTF.Shared.MUTEX.EBH\n  ShimCacheMutex\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n\nMUTEXES OPENED (2):\n  \\Sessions\\1\\BaseNamedObjects\\Global\\TermSrvReadyEvent\n  \\Sessions\\1\\BaseNamedObjects\\Local\\MSCTF.CtfActivated.Default1\n\nPROCESSES CREATED (17):\n  %WINDIR%\\explorer.exe\n  %WINDIR%\\syswow64\\cmd.exe\n  %WINDIR%\\syswow64\\explorer.exe\n  %WINDIR%\\syswow64\\reg.exe\n  %WINDIR%\\syswow64\\xcopy.exe\n  <PATH_SAMPLE.EXE>\n  <SYSTEM32>\\conhost.exe\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\3.tmp\\promo.bat\" \"\n  C:\\WINDOWS\\explorer.exe\n  C:\\WINDOWS\\explorer.exe explorer  Promo\\\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress\n\nCOMPILER AND PACKER:\n  [Compiler] PureBasic\n\nPE SECTIONS (entropy):\n  .code        entropy=5.56\n  .text        entropy=6.55\n  .rdata       entropy=7.28\n  .data        entropy=4.91\n  .rsrc        entropy=6.79\n",
  "ground_truth": {
    "hash": "3657e8733cdea568bc631f2df3ec1bae1858e31e5c88cabb84558ce36bc42094",
    "verdict": "MALICIOUS",
    "detection_rate": 0.7397,
    "malicious_count": 54,
    "popular_threat_category": [
      {
        "value": "trojan",
        "count": 13
      },
      {
        "value": "miner",
        "count": 2
      },
      {
        "value": "virus",
        "count": 2
      }
    ],
    "popular_threat_name": [
      {
        "value": "virtob",
        "count": 7
      },
      {
        "value": "virut",
        "count": 5
      },
      {
        "value": "hesv",
        "count": 4
      }
    ],
    "suggested_label": "trojan.virtob/virut",
    "av_family_vote": "win32",
    "has_family_label": true,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via polling"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "check OS version"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1497",
        "signature_description": "check for unmoving mouse cursor"
      },
      {
        "id": "T1547",
        "signature_description": "reference startup folder"
      },
      {
        "id": "T1614",
        "signature_description": "identify system language via API"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create process on Windows",
      "delay execution",
      "create or open file",
      "encode data using XOR",
      "allocate thread local storage",
      "contain loop",
      "identify system language via API",
      "delete directory",
      "get common file path",
      "extract resource via kernel32 functions",
      "link function at runtime on Windows",
      "enumerate gui resources",
      "accept command line arguments",
      "terminate process",
      "get OS version",
      "check OS version",
      "get graphical window text",
      "log keystrokes via polling",
      "check for unmoving mouse cursor",
      "delete file",
      "create directory",
      "write file on Windows",
      "read file on Windows",
      "hash data with CRC32",
      "hash data with MD5",
      "contain a resource (.rsrc) section",
      "reference startup folder",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "<HKLM>\\SOFTWARE\\Wow6432Node\\Microsoft\\Windows\\CurrentVersion\\Run\\Promo",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Promo",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\Promo"
      ],
      "registry_all_count": 6,
      "files_dropped": [
        "%TEMP%\\61bb.tmp\\promo.bat"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "MSCTF.Shared.MUTEX.EBH",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "%WINDIR%\\explorer.exe",
        "%WINDIR%\\syswow64\\cmd.exe",
        "%WINDIR%\\syswow64\\explorer.exe",
        "%WINDIR%\\syswow64\\reg.exe",
        "%WINDIR%\\syswow64\\xcopy.exe",
        "<PATH_SAMPLE.EXE>",
        "<SYSTEM32>\\conhost.exe",
        "C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\3.tmp\\promo.bat\" \"",
        "C:\\WINDOWS\\explorer.exe",
        "C:\\WINDOWS\\explorer.exe explorer  Promo\\",
        "C:\\WINDOWS\\system32\\cmd.exe",
        "C:\\WINDOWS\\system32\\cmd.exe C:\\WINDOWS\\system32\\cmd.exe /c echo 12:42:54.54",
        "C:\\WINDOWS\\system32\\cmd.exe C:\\WINDOWS\\system32\\cmd.exe /c echo 12:42:54.58",
        "C:\\WINDOWS\\system32\\cmd.exe C:\\WINDOWS\\system32\\cmd.exe /c echo Fri 11/15/2019",
        "C:\\WINDOWS\\system32\\cmd.exe C:\\WINDOWS\\system32\\cmd.exe /c type C:\\Documents and Settings\\<USER>\\Promo_time\"\"",
        "C:\\WINDOWS\\system32\\reg.exe",
        "C:\\WINDOWS\\system32\\xcopy.exe"
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Dr.Web vxCube",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}