# dataset/preprocessed/36f654553778978b0e83342fe6ce6019cff581a6865ceb97adee7c0eb1b80dbc.json
{
  "hash": "36f654553778978b0e83342fe6ce6019cff581a6865ceb97adee7c0eb1b80dbc",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\36f654553778978b0e83342fe6ce6019cff581a6865ceb97adee7c0eb1b80dbc(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\36f654553778978b_chinese cumshot blowjob girls hotel .zip.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american fetish blowjob uncut traffic .mpeg.exe  [sha256: 71b93e1b1be8ac2b...]\n  black gang bang xxx [free] feet .rar.exe  [sha256: eb715854d939a462...]\n  black porn fucking girls cock .mpeg.exe  [sha256: 1ede1a4d0e54a1bd...]\n  brasilian handjob xxx masturbation sweet .rar.exe  [sha256: 2ccec00f2f5cf3db...]\n  brasilian porn trambling [free]  (Liz).avi.exe  [sha256: 1d1a575482579fbf...]\n  danish nude bukkake masturbation shoes .avi.exe  [sha256: 8dff7974ec947c4c...]\n  defrag64.exe  [sha256: 82aa9779ac98b42a...]\n  deinstall.exe  [sha256: 2d17d959e4e1c07c...]\n  driver.exe  [sha256: b1f26f6eb309d7e4...]\n  drwatson32.exe  [sha256: cf38e8c533d4b983...]\n  gay catfight redhair .mpeg.exe  [sha256: 5c9399dd9f66778a...]\n  hardcore girls boots .mpeg.exe  [sha256: 529a576e63633a3f...]\n  horse hidden .mpeg.exe  [sha256: 2d6f308bf5f27244...]\n  IE9setup.exe  [sha256: c44b400d74e90088...]\n  italian cumshot gay catfight hole wifey  (Janette).mpg.exe  [sha256: 148e0fd19c222456...]\n  italian handjob xxx full movie hairy  (Sandy,Jade).rar.exe  [sha256: ce9dd611a61ed492...]\n  japanese action bukkake masturbation .mpeg.exe  [sha256: f0f8a954a0924d21...]\n  japanese action horse lesbian cock femdom  (Sarah).zip.exe  [sha256: 7276baa9b5670ecd...]\n  japanese action lingerie uncut hairy .avi.exe  [sha256: 12bd99631407d7ba...]\n  japanese beastiality lingerie hidden glans beautyfull .rar.exe  [sha256: fdd243e477fc17e0...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\japanese fetish fucking uncut sweet .rar.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\swedish horse trambling [bangbus] .mpg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\xxx public hole  (Sandy,Janette).avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\danish nude bukkake masturbation shoes .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\japanese cum horse masturbation circumcision .rar.exe\n  C:\\Documents and Settings\\Default User\\Templates\\horse hidden .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\japanese action lingerie uncut hairy .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\trambling catfight beautyfull .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\russian animal lesbian public boots .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\hardcore girls boots .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\russian fetish horse big titts .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\black gang bang xxx [free] feet .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\tyrkish porn blowjob catfight .zip.exe\n  C:\\WINDOWS\\Temp\\italian cumshot gay catfight hole wifey  (Janette).mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\black porn fucking girls cock .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\brasilian handjob xxx masturbation sweet .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\japanese action bukkake masturbation .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\gay catfight redhair .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\japanese action horse lesbian cock femdom  (Sarah).zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.879939079284668  size=1249378\n",
  "ground_truth": {
    "hash": "36f654553778978b0e83342fe6ce6019cff581a6865ceb97adee7c0eb1b80dbc",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9571,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\american blowjob kicking voyeur .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\american fetish [free] .mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\norwegian action lingerie masturbation feet .zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\malaysia cum kicking masturbation nipples hairy  (Britney).rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\spanish horse catfight .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\hardcore beast full movie boobs 50+ .rar.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american fetish blowjob uncut traffic .mpeg.exe",
        "black gang bang xxx [free] feet .rar.exe",
        "black porn fucking girls cock .mpeg.exe",
        "brasilian handjob xxx masturbation sweet .rar.exe",
        "brasilian porn trambling [free]  (Liz).avi.exe",
        "danish nude bukkake masturbation shoes .avi.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe",
        "gay catfight redhair .mpeg.exe",
        "hardcore girls boots .mpeg.exe",
        "horse hidden .mpeg.exe",
        "italian cumshot gay catfight hole wifey  (Janette).mpg.exe",
        "italian handjob xxx full movie hairy  (Sandy,Jade).rar.exe",
        "japanese action bukkake masturbation .mpeg.exe",
        "japanese action horse lesbian cock femdom  (Sarah).zip.exe",
        "japanese action lingerie uncut hairy .avi.exe",
        "japanese beastiality lingerie hidden glans beautyfull .rar.exe",
        "japanese cum horse masturbation circumcision .rar.exe",
        "japanese fetish fucking uncut sweet .rar.exe",
        "japanese nude bukkake public sweet .mpeg.exe",
        "lingerie uncut penetration .rar.exe",
        "mssrv.exe",
        "rundll64.exe",
        "russian animal lesbian public boots .mpeg.exe",
        "russian fetish horse big titts .zip.exe",
        "russian kicking xxx big  (Janette).rar.exe",
        "setup.exe",
        "soccer.exe",
        "start.exe",
        "swedish horse trambling [bangbus] .mpg.exe",
        "trambling catfight beautyfull .avi.exe",
        "trambling masturbation  (Samantha).rar.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}