# dataset/preprocessed/36ff3b26618cc6fa1d4c1259c77cb734b2809fdeff82c46e43d2e484a2c633ca.json
{
  "hash": "36ff3b26618cc6fa1d4c1259c77cb734b2809fdeff82c46e43d2e484a2c633ca",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n    C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n    C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n      C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n    C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n      C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\36ff3b26618cc6fa1d4c1259c77cb734b2809fdeff82c46e43d2e484a2c633ca(null) [(null)]\n  C:\\Program Files\\dotnet\\shared\\swedish horse [bangbus] fishy  (Janette,Kathrin).rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american trambling fetish sleeping hole .avi.exe  [sha256: 245efe697992974e...]\n  black animal sleeping sweet  (Sylvia).mpg.exe  [sha256: 7bbce49187baa734...]\n  brasilian cumshot nude full movie blondie .mpeg.exe  [sha256: 27a225d828b902ee...]\n  brasilian xxx action hidden vagina wifey .mpeg.exe  [sha256: b971503eaedf0e9f...]\n  british sperm kicking [bangbus] titts lady .mpeg.exe  [sha256: 2c8270b6c453904b...]\n  bukkake hardcore girls penetration .avi.exe  [sha256: 5b8dfff9e0ef8c94...]\n  bukkake masturbation 50+ .avi.exe  [sha256: dd760be993928d0f...]\n  canadian hardcore full movie legs .rar.exe  [sha256: fec55bc498919736...]\n  chinese lesbian masturbation sweet  (Sandy).rar.exe  [sha256: 82b5727685451e20...]\n  defrag.exe  [sha256: f1008ad08a548302...]\n  deinstall.exe  [sha256: 3992b3645af52718...]\n  drwatson32.exe  [sha256: d4b3fdfb47485486...]\n  edit.exe  [sha256: 2a754e9325275b9a...]\n  fetish uncut fishy .zip.exe  [sha256: 472865c147c65ac5...]\n  game3d.exe  [sha256: aee5a45b06a52643...]\n  gang bang beastiality licking nipples  (Curtney,Samantha).mpg.exe  [sha256: cf131028345b87c6...]\n  gang bang lesbian .rar.exe  [sha256: 56cac62a5af224dd...]\n  german lesbian gay girls beautyfull .mpg.exe  [sha256: ee1514ea504cd2e1...]\n  handjob fetish sleeping stockings .zip.exe  [sha256: 87ba9e681c181647...]\n  hardcore big redhair  (Gina,Janette).zip.exe  [sha256: 0054bf5f7aabe6e2...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Templates\\chinese lesbian masturbation sweet  (Sandy).rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\german lesbian gay girls beautyfull .mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\swedish beastiality lesbian [milf] .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\gang bang beastiality licking nipples  (Curtney,Samantha).mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\spanish gay bukkake full movie .mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\bukkake masturbation 50+ .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\american trambling fetish sleeping hole .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\gang bang lesbian .rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\brasilian xxx action hidden vagina wifey .mpeg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\bukkake hardcore girls penetration .avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\indian porn catfight gorgeoushorny .mpg.exe\n  C:\\WINDOWS\\Temp\\hardcore big redhair  (Gina,Janette).zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\tyrkish action horse girls .avi.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\tyrkish horse [milf] feet .mpg.exe\n  C:\\WINDOWS\\assembly\\tmp\\handjob fetish sleeping stockings .zip.exe\n  C:\\WINDOWS\\ime\\shared\\canadian hardcore full movie legs .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\porn sleeping boobs mistress .mpeg.exe\n  C:\\WINDOWS\\security\\templates\\british sperm kicking [bangbus] titts lady .mpeg.exe\n  C:\\WINDOWS\\system32\\config\\systemprofile\\Local Settings\\Temp\\horse masturbation black hairunshaved .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.88010311126709  size=774933\n",
  "ground_truth": {
    "hash": "36ff3b26618cc6fa1d4c1259c77cb734b2809fdeff82c46e43d2e484a2c633ca",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\german sperm masturbation .mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\beast bukkake lesbian hotel  (Melissa,Jade).mpeg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\chinese xxx uncut feet mistress  (Sandy).rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\british beastiality girls mature .rar.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\kicking big high heels .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\chinese porn beast uncut young  (Janette).zip.exe",
        "C:\\Windows\\mssrv.exe",
        "american trambling fetish sleeping hole .avi.exe",
        "black animal sleeping sweet  (Sylvia).mpg.exe",
        "brasilian cumshot nude full movie blondie .mpeg.exe",
        "brasilian xxx action hidden vagina wifey .mpeg.exe",
        "british sperm kicking [bangbus] titts lady .mpeg.exe",
        "bukkake hardcore girls penetration .avi.exe",
        "bukkake masturbation 50+ .avi.exe",
        "canadian hardcore full movie legs .rar.exe",
        "chinese lesbian masturbation sweet  (Sandy).rar.exe",
        "defrag.exe",
        "deinstall.exe",
        "drwatson32.exe",
        "edit.exe",
        "fetish uncut fishy .zip.exe",
        "game3d.exe",
        "gang bang beastiality licking nipples  (Curtney,Samantha).mpg.exe",
        "gang bang lesbian .rar.exe",
        "german lesbian gay girls beautyfull .mpg.exe",
        "handjob fetish sleeping stockings .zip.exe",
        "hardcore big redhair  (Gina,Janette).zip.exe",
        "horse masturbation black hairunshaved .rar.exe",
        "indian porn catfight gorgeoushorny .mpg.exe",
        "malaysia sperm big .avi.exe",
        "mssrv.exe",
        "nude beast [milf] 50+ .rar.exe",
        "porn sleeping boobs mistress .mpeg.exe",
        "safetyserver.exe",
        "screensaver.scr",
        "soccer.exe",
        "spanish gay bukkake full movie .mpg.exe",
        "speedup.exe",
        "summergames.exe",
        "swedish animal lesbian catfight ash ejaculation .mpg.exe",
        "swedish beastiality lesbian [milf] .mpeg.exe",
        "tyrkish action horse girls .avi.exe",
        "tyrkish beast cum [milf] legs 50+  (Jenna).avi.exe",
        "tyrkish horse [milf] feet .mpg.exe",
        "undelete.exe",
        "unformat.exe",
        "word.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}