# dataset/preprocessed/39475a70c86b5063fea0350ff71232f8615de107ce15fb93d43ef49d63e8e311.json
{
  "hash": "39475a70c86b5063fea0350ff71232f8615de107ce15fb93d43ef49d63e8e311",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  american cum lingerie [milf] gorgeoushorny .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\american cum lingerie [milf] gorgeoushorny .mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\american cum lingerie [milf] gorgeoushorny .mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\39475a70c86b5063fea0350ff71232f8615de107ce15fb93d43ef49d63e8e311(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\american cum lingerie [milf] gorgeoushorny .mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american fetish fucking lesbian sweet .zip.exe  [sha256: fedbaf2d8718109c...]\n  american gang bang hardcore big cock .rar.exe  [sha256: 870ed6dc44cb92cf...]\n  american horse gay lesbian  (Curtney).mpg.exe  [sha256: 5d608d92198766ad...]\n  american kicking beast licking .avi.exe  [sha256: dc8e15d7f687c4b1...]\n  american kicking gay catfight 50+ .avi.exe  [sha256: ea46e407b99e05d4...]\n  american porn xxx uncut  (Karin).zip.exe  [sha256: aed27439c24ff9ae...]\n  beast licking feet upskirt  (Janette).mpeg.exe  [sha256: 3109417669d351ac...]\n  black kicking lingerie catfight blondie .mpeg.exe  [sha256: 6e7fda1376d2f08c...]\n  blowjob public  (Jade).mpg.exe  [sha256: 1129b146a880647f...]\n  brasilian cumshot fucking [bangbus]  (Jade).mpg.exe  [sha256: 8274d64240829b60...]\n  brasilian fetish xxx hot (!) titts sweet  (Samantha).mpg.exe  [sha256: 8bf514fc90096879...]\n  brasilian horse trambling catfight cock .zip.exe  [sha256: a9a069c2749adee3...]\n  danish beastiality horse voyeur titts .mpeg.exe  [sha256: 127f3cff6eb0e49c...]\n  danish beastiality lingerie girls wifey .rar.exe  [sha256: 5f54d1a3b06d4799...]\n  danish kicking lingerie lesbian .zip.exe  [sha256: ec4a4bdc4db4d94a...]\n  defrag.exe  [sha256: a648614305a285e9...]\n  defrag64.exe  [sha256: e28e555638ab73d7...]\n  deinstall.exe  [sha256: 8fcb2ab663344746...]\n  driver.exe  [sha256: 73bff1e44fa31f9b...]\n  fucking big latex  (Sonja,Melissa).zip.exe  [sha256: 204a605ceb6c17fe...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\danish beastiality horse voyeur titts .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\russian fetish gay [milf] .rar.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\japanese handjob xxx hidden cock traffic .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\italian handjob gay lesbian mature .rar.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\horse girls ash  (Jenna,Jade).mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\american kicking beast licking .avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\indian handjob lingerie [bangbus] feet young  (Janette).mpg.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\blowjob public  (Jade).mpg.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\beast licking feet upskirt  (Janette).mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\american fetish fucking lesbian sweet .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\american gang bang hardcore big cock .rar.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\swedish gang bang gay big glans high heels .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\danish kicking lingerie lesbian .zip.exe\n  C:\\WINDOWS\\Temp\\italian kicking lesbian public blondie .mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\brasilian fetish xxx hot (!) titts sweet  (Samantha).mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\tyrkish action fucking public glans .rar.exe\n  C:\\WINDOWS\\assembly\\temp\\american kicking gay catfight 50+ .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\american porn xxx uncut  (Karin).zip.exe\n  C:\\WINDOWS\\ime\\shared\\fucking big latex  (Sonja,Melissa).zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.879557132720947  size=1921616\n",
  "ground_truth": {
    "hash": "39475a70c86b5063fea0350ff71232f8615de107ce15fb93d43ef49d63e8e311",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\sperm trambling [milf] .avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\malaysia cum cum voyeur hotel .zip.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\american horse masturbation vagina .zip.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\nude uncut ash mature .avi.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\horse uncut vagina stockings  (Sonja,Sarah).mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\russian xxx nude full movie .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\cum [bangbus] ash 50+ .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\gang bang catfight hotel .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\template\\russian trambling cumshot big legs swallow .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\xxx several models  (Jade).mpg.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\swedish cumshot fetish several models .zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\norwegian sperm [bangbus] .mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\blowjob catfight swallow  (Melissa,Janette).zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\xxx uncut boobs ash .avi.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\norwegian fetish kicking [milf] balls .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\swedish cumshot voyeur redhair  (Janette).rar.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american fetish fucking lesbian sweet .zip.exe",
        "american gang bang hardcore big cock .rar.exe",
        "american horse gay lesbian  (Curtney).mpg.exe",
        "american kicking beast licking .avi.exe",
        "american kicking gay catfight 50+ .avi.exe",
        "american porn xxx uncut  (Karin).zip.exe",
        "beast licking feet upskirt  (Janette).mpeg.exe",
        "black kicking lingerie catfight blondie .mpeg.exe",
        "blowjob public  (Jade).mpg.exe",
        "brasilian cumshot fucking [bangbus]  (Jade).mpg.exe",
        "brasilian fetish xxx hot (!) titts sweet  (Samantha).mpg.exe",
        "brasilian horse trambling catfight cock .zip.exe",
        "danish beastiality horse voyeur titts .mpeg.exe",
        "danish beastiality lingerie girls wifey .rar.exe",
        "danish kicking lingerie lesbian .zip.exe",
        "defrag.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "fucking big latex  (Sonja,Melissa).zip.exe",
        "gay [free] cock high heels  (Samantha).mpg.exe",
        "horse [bangbus] hairy .mpg.exe",
        "horse girls ash  (Jenna,Jade).mpg.exe",
        "indian animal fucking full movie cock ash .zip.exe",
        "indian handjob lingerie [bangbus] feet young  (Janette).mpg.exe",
        "italian handjob gay lesbian mature .rar.exe",
        "italian kicking lesbian public blondie .mpg.exe",
        "japanese cum sperm hidden .zip.exe",
        "japanese handjob xxx hidden cock traffic .mpeg.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "russian fetish gay [milf] .rar.exe",
        "screensaver.scr",
        "soccer.exe",
        "swedish gang bang gay big glans high heels .zip.exe",
        "swedish horse horse [bangbus]  (Karin).mpg.exe",
        "tyrkish action fucking public glans .rar.exe",
        "undelete.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}