# dataset/preprocessed/397dd651360321c6c62976b236135319d62f18e7d643413d4dd25962b7490bb7.json
{
  "hash": "397dd651360321c6c62976b236135319d62f18e7d643413d4dd25962b7490bb7",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  ****.exe\n    tqbjx.exe\n\nCOMMAND EXECUTIONS (1):\n  (open) C:\\397dd651360321c6c62976b236135319d62f18e7d643413d4dd25962b74(null) [(null)]\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (6):\n  - contain obfuscated stackstrings\n  - get startup folder\n  - set file attributes\n  - get common file path\n  - accept command line arguments\n  - packed with generic packer\n\nCAPA RULE MATCHES (MBC):\n  - Anti-Static Analysis::Executable Code Obfuscation::Argument Obfuscation [B0032.020]\n  - Anti-Static Analysis::Executable Code Obfuscation::Stack Strings [B0032.017]\n  - get startup folder\n  - Process::Create Process [C0017]\n  - File System::Set File Attributes [C0050]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - contain loop\n  - Discovery::File and Directory Discovery [E1083]\n  - File System::Delete File [C0047]\n  - File System::Create Directory [C0046]\n  - Process::Allocate Thread Local Storage [C0040]\n  - get thread local storage value\n  - Process::Set Thread Local Storage Value [C0041]\n  - File System::Read File [C0051]\n  - Execution::Command and Scripting Interpreter [E1059]\n  - contain pusha popa sequence\n  - resolve function by parsing PE exports\n  - Anti-Static Analysis::Software Packing::Standard Compression [F0001.002]\n  - contain a resource (.rsrc) section\n  - (internal) packer file limitation\n\n--- ARTIFACTS ---\n\nFILES DROPPED (23):\n  397dd651360321c6c62976b236135319d62f18e7d643413d4dd25962b7490bb7 .exe  [sha256: 9bc18a1b4689756d...]\n  AUTOEXEC.BAT .exe  [sha256: 5a7212be51691252...]\n  boot.ini .exe  [sha256: 577db3875b1b39d3...]\n  CONFIG.SYS .exe  [sha256: ab7a9f8af2144c7e...]\n  cuckoo .exe  [sha256: e043dd5ca64ca5dc...]\n  Documents and Settings .exe  [sha256: f3ae798af9ed18a2...]\n  eula.1028.txt .exe  [sha256: c088d644b2b21340...]\n  eula.1031.txt .exe  [sha256: 595f4692e8f4aba8...]\n  eula.1033.txt .exe  [sha256: c6647c287ef314e6...]\n  eula.1036.txt .exe  [sha256: 32b80bf9a7605fb7...]\n  eula.1040.txt .exe  [sha256: 67ac7a44fa76e047...]\n  eula.1041.txt .exe  [sha256: 6c46ea43eaaf5d09...]\n  eula.1042.txt .exe  [sha256: f2db974106c12e08...]\n  eula.2052.txt .exe  [sha256: adf4dc33fda48233...]\n  eula.3082.txt .exe  [sha256: c4ef5c25806146b0...]\n  globdata.ini .exe  [sha256: 693c351d3013ba18...]\n  install.exe .exe  [sha256: bf48eef03ad69aca...]\n  install.ini .exe  [sha256: 982f666f72080b3b...]\n  install.res.1028.dll .exe  [sha256: 66e8528c508233ed...]\n  install.res.1031.dll .exe  [sha256: f933f43c74ab2ee9...]\n\nFILES WRITTEN (30):\n  C:\\0 .exe\n  C:\\AUTOEXEC.BAT .exe\n  C:\\CONFIG.SYS .exe\n  C:\\DiskD .exe\n  C:\\DiskX .exe\n  C:\\DllLoader.exe .exe\n  C:\\Documents and Settings .exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Saaaalamm\\Mira.h\n  C:\\Documents and Settings\\All Users\\Application Data\\tqbjx.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\xamish.exe\n  C:\\IO.SYS .exe\n  C:\\MSDOS.SYS .exe\n  C:\\Mirav-newfile .exe\n  C:\\Mirav-samplefile .exe\n  C:\\NTDETECT.COM .exe\n  C:\\Program Files .exe\n  C:\\boot.ini .exe\n  C:\\bootfont.bin .exe\n  C:\\eula.2052.txt .exe\n  C:\\globdata.ini .exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System = C:\\Documents and Settings\\All Users\\Application Data\\tqbjx.exe\u0000\n\nMUTEXES CREATED (7):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  C:\\Documents and Settings\\All Users\\Application Data\\xamish.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA\n\nCOMPILER AND PACKER:\n  [Compiler] MinGW\n  [Linker] GNU linker ld (GNU Binutils)\n\nPE SECTIONS (entropy):\n  .text        entropy=6.11\n  .data        entropy=0.71\n  .rdata       entropy=5.01\n  .bss         entropy=0.0\n  .idata       entropy=4.29\n  .rsrc        entropy=6.21\n  [overlay]    entropy=7.022463798522949  size=264655\n",
  "ground_truth": {
    "hash": "397dd651360321c6c62976b236135319d62f18e7d643413d4dd25962b7490bb7",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8611,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "contain obfuscated stackstrings"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "get startup folder"
      }
    ],
    "capa_sig_names": [
      "contain obfuscated stackstrings",
      "get startup folder",
      "create process on Windows",
      "set file attributes",
      "delay execution",
      "create or open file",
      "contain loop",
      "get common file path",
      "delete file",
      "create directory",
      "allocate thread local storage",
      "get thread local storage value",
      "set thread local storage value",
      "read file on Windows",
      "accept command line arguments",
      "contain pusha popa sequence",
      "resolve function by parsing PE exports",
      "packed with generic packer",
      "contain a resource (.rsrc) section",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "\\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Microsoft?Windows?Operating System"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "397dd651360321c6c62976b236135319d62f18e7d643413d4dd25962b7490bb7 .exe",
        "AUTOEXEC.BAT .exe",
        "CONFIG.SYS .exe",
        "Documents and Settings .exe",
        "boot.ini .exe",
        "cuckoo .exe",
        "eula.1028.txt .exe",
        "eula.1031.txt .exe",
        "eula.1033.txt .exe",
        "eula.1036.txt .exe",
        "eula.1040.txt .exe",
        "eula.1041.txt .exe",
        "eula.1042.txt .exe",
        "eula.2052.txt .exe",
        "eula.3082.txt .exe",
        "globdata.ini .exe",
        "install.exe .exe",
        "install.ini .exe",
        "install.res.1028.dll .exe",
        "install.res.1031.dll .exe",
        "install.res.1033.dll .exe",
        "install.res.1036.dll .exe",
        "install.res.1040.dll .exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\Documents and Settings\\All Users\\Application Data\\xamish.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork"
  ]
}