# dataset/preprocessed/39b966ed5ce72ff5fe3bce967854bead60bde310eeb60cac00cabf28bbd04872.json
{
  "hash": "39b966ed5ce72ff5fe3bce967854bead60bde310eeb60cac00cabf28bbd04872",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (3):\n  ftp.myfreesite.foreverhost.us\n  ftp5.nofeehost.com\n  ftp5.nofeehost.com \u2192 35.186.238.101\n\nSERVICES OPENED (1):\n  RASMAN\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  <PATH_SAMPLE.EXE>\n    %ProgramFiles(x86)%\\3bcdc871\\jusched.exe\n  <SYSTEM32>\\svchost.exe\n  %ProgramFiles(x86)%\\3bcdc871\\jusched.exe\n\nCOMMAND EXECUTIONS (2):\n  (open) C:\\Program Files\\6c78a9fa\\jusched.exe(null) [(null)]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (5):\n  - get geographical location\n  - enumerate files on Windows\n  - parse PE header\n  - schedule task via ITaskScheduler\n  - get session user name\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - contain loop\n  - get geographical location\n  - Discovery::File and Directory Discovery [E1083]\n  - parse PE header\n  - resolve function by parsing PE exports\n  - schedule task via ITaskScheduler\n  - get session user name\n  - debug build\n  - contain a resource (.rsrc) section\n\n--- ARTIFACTS ---\n\nFILES DROPPED (7):\n  jusched.exe  [sha256: 5ab6712024f1eec2...]\n  C:\\Program Files (x86)\\71d2c87e\\jusched.exe  [sha256: ee39b096abadf8b2...]\n  %ProgramFiles(x86)%\\3bcdc871\\3bcdc871  [sha256: 6bc8a4f93eaa1b3e...]\n  %ProgramFiles(x86)%\\3bcdc871\\info_a  [sha256: 5cb4b814c1dd88a2...]\n  %ProgramFiles(x86)%\\3bcdc871\\jusched.exe  [sha256: 85b358c86957eaf7...]\n  <SYSTEM32>\\tasks\\update23  [sha256: 6bdbcf8f2582738c...]\n  %WINDIR%\\tasks\\update23.job  [sha256: e59b545a468a59ec...]\n\nFILES WRITTEN (10):\n  C:\\Program Files (x86)\\71d2c87e\\71d2c87e\n  C:\\Program Files (x86)\\71d2c87e\\info_a\n  C:\\Program Files\\6c78a9fa\\6c78a9fa\n  C:\\Program Files\\6c78a9fa\\info_a\n  C:\\Program Files\\6c78a9fa\\jusched.exe\n  C:\\Program Files\\f440491e\\f440491e\n  C:\\Program Files\\f440491e\\info_a\n  C:\\Program Files\\f440491e\\jusched.exe\n  C:\\WINDOWS\\Tasks\\Update23.job\n  C:\\Windows\\Tasks\\Update23.job\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe = C:\\Program Files (x86)\\71d2c87e\\jusched.exe:*:Enabled:JavaUpdate23\n  \\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe = C:\\Program Files\\f440491e\\jusched.exe:*:Enabled:JavaUpdate23\u0000\n\nOTHER REGISTRY KEYS SET (10):\n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\ShellNoRoam\\MUICache\\C:\\Program \n  \\REGISTRY\\USER\\S-1-5-21-1482476501-1645522239-1417001333-500\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\UNCAsIntranet\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\AutoDetect\n  \\REGISTRY\\A\\{61610fc8-fa78-11e9-8fef-cc1f8471db16}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectId_\n  \\REGISTRY\\A\\{61610fc8-fa78-11e9-8fef-cc1f8471db16}\\DefaultObjectStore\\LruList\\CurrentLru\n  \\REGISTRY\\A\\{61610fc8-fa78-11e9-8fef-cc1f8471db16}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectId\n  \\REGISTRY\\A\\{61610fc8-fa78-11e9-8fef-cc1f8471db16}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectLru\n  \\REGISTRY\\A\\{61610fc8-fa78-11e9-8fef-cc1f8471db16}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectLru_\n  \\REGISTRY\\A\\{61610fc8-fa78-11e9-8fef-cc1f8471db16}\\DefaultObjectStore\\ObjectTable\\E9\\_FileId_\n\nMUTEXES CREATED (11):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  RasPbFile\n  ShimCacheMutex\n  ZonesCacheCounterMutex\n  ZonesCounterMutex\n\nMUTEXES OPENED (4):\n  _!MSFTHISTORY!_\n  c:!documents and settings!administrator!cookies!\n  c:!documents and settings!administrator!local settings!history!history.ie5!\n  c:!documents and settings!administrator!local settings!temporary internet files!content.ie5!\n\nPROCESSES CREATED (3):\n  %ProgramFiles(x86)%\\3bcdc871\\jusched.exe\n  <PATH_SAMPLE.EXE>\n  C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \"\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  .text        entropy=3.92\n  .rdata       entropy=4.65\n  .data        entropy=0.87\n  .idata       entropy=3.03\n  .rsrc        entropy=3.08\n  .reloc       entropy=1.21\n  myvzbsb      entropy=0.03\n  yllgqsl      entropy=0.05\n  asciors      entropy=0.0\n  .drdata      entropy=0.0\n  [overlay]    entropy=6.9600043296813965  size=4508\n",
  "ground_truth": {
    "hash": "39b966ed5ce72ff5fe3bce967854bead60bde310eeb60cac00cabf28bbd04872",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9155,
    "malicious_count": 65,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via ITaskScheduler"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1129",
        "signature_description": "parse PE header"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "contain loop",
      "get geographical location",
      "enumerate files on Windows",
      "parse PE header",
      "resolve function by parsing PE exports",
      "schedule task via ITaskScheduler",
      "get session user name",
      "debug build",
      "contain a resource (.rsrc) section"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "ftp.myfreesite.foreverhost.us",
        "ftp5.nofeehost.com"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SYSTEM\\CurrentControlSet\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "\\REGISTRY\\MACHINE\\SYSTEM\\ControlSet002\\Services\\SharedAccess\\Parameters\\FirewallPolicy\\StandardProfile\\AuthorizedApplications\\List\\C:\\Program Files\\f440491e\\jusched.exe"
      ],
      "registry_all_count": 37,
      "files_dropped": [
        "%ProgramFiles(x86)%\\3bcdc871\\3bcdc871",
        "%ProgramFiles(x86)%\\3bcdc871\\info_a",
        "%ProgramFiles(x86)%\\3bcdc871\\jusched.exe",
        "%WINDIR%\\tasks\\update23.job",
        "<SYSTEM32>\\tasks\\update23",
        "C:\\Program Files (x86)\\71d2c87e\\jusched.exe",
        "jusched.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "RasPbFile",
        "ShimCacheMutex",
        "ZonesCacheCounterMutex",
        "ZonesCounterMutex",
        "ZonesLockedCacheCounterMutex"
      ],
      "processes": [
        "%ProgramFiles(x86)%\\3bcdc871\\jusched.exe",
        "<PATH_SAMPLE.EXE>",
        "C:\\Program Files\\6c78a9fa\\jusched.exe C:\\Program Files\\6c78a9fa\\jusched.exe\" \""
      ]
    }
  },
  "num_sandboxes": 5,
  "sandbox_names": [
    "CAPA",
    "Dr.Web vxCube",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}