# dataset/preprocessed/3a863472603522f8123dbc92d1cfc2f606aae2ef61c66623c05856b02dda5425.json
{
  "hash": "3a863472603522f8123dbc92d1cfc2f606aae2ef61c66623c05856b02dda5425",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  232.csv.exe\n    C:\\Windows\\system32\\HelpMe.exe\n\nCOMMAND EXECUTIONS (1):\n  C:\\Windows\\system32\\HelpMe.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (2):\n  - packed with ASPack\n  - reference startup folder\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - contain loop\n  - contain a thread local storage (.tls) section\n  - Anti-Static Analysis::Software Packing [F0001]\n  - contain a resource (.rsrc) section\n  - reference startup folder\n  - (internal) packer file limitation\n\n--- ARTIFACTS ---\n\nFILES DROPPED (17):\n  420.csv.exe  [sha256: 001ef4cc95a83b91...]\n  664.csv.exe  [sha256: 07c768f6e84347fb...]\n  AUTORUN.INF.exe  [sha256: 4ad9ae99bf25166c...]\n  chrome.exe.exe  [sha256: e29683c0d07e4eae...]\n  DefaultStore_59R.bin.exe  [sha256: fc9f1bf630281539...]\n  DUlHWC.dll.exe  [sha256: 1363e77beceb8c97...]\n  HelpMe.exe  [sha256: 389274da044612eb...]\n  index.dat.exe  [sha256: 67e74269045afec8...]\n  UserMigratedStore_59R.bin.exe  [sha256: 6126a5d543bc8038...]\n  wordpfct.wpd.exe  [sha256: a99d706fb279507d...]\n  wordpfct.wpg.exe  [sha256: a8ee608ef6a94962...]\n  C:\\AUTORUN.INF.exe  [sha256: bb8574c9f077cb89...]\n  \\AUTORUN.INF  [sha256: cfe448b4506a95b3...]\n  \\MSOCache\\All Users\\{90140000-001B-0409-0000-0000000FF1CE}-C\\WordMUI.xml.exe  [sha256: 3cd1accfe9674f57...]\n  \\AutoRun.exe  [sha256: 3a863472603522f8...]\n  \\AUTORUN.INF.exe  [sha256: 61b2ce95cdf00c14...]\n  \\Windows\\System32\\HelpMe.exe  [sha256: 389274da044612eb...]\n\nFILES WRITTEN (30):\n  C:\\AUTOEXEC.BAT.exe\n  C:\\AUTORUN.INF.exe\n  C:\\CONFIG.SYS.exe\n  C:\\Documents and Settings\\<USER>\\Start Menu\\Programs\\Startup\\Soft.lnk\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Crypto\\RSA\\S-1-5-18\\d42cc0c3858a58db2db37658219e6400_fabbc6a1-c573-4ea0-9ca1-50004b35a440.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Media Player\\DefaultStore_59R.bin.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Media Player\\UserMigratedStore_59R.bin.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Connections\\Pbk\\sharedaccess.ini.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\qmgr0.dat.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\qmgr1.dat.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\<USER>.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\Default Pictures\\airplane.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\Default Pictures\\astronaut.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\guest.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\Beethoven's Symphony No. 9 (Scherzo).wma.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\New Stories (Highway Blues).wma.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Playlists\\000511FF\\Plylst1.wpl.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SoftWare\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell = Explorer.exe  HelpMe.exe\n\nOTHER REGISTRY KEYS SET (3):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Explorer\\Advanced\\Folder\\Hidden\\SHOWALL\\CheckedValue\n  \\REGISTRY\\MACHINE\\HexStream = 53004F004600540057004100520045005C004D006900630072006F0073006F00660074005C00570069006E0064\n  \\REGISTRY\\MACHINE\\HexStream = 53006F006600740057006100720065005C004D006900630072006F0073006F00660074005C00570069006E0064\n\nMUTEXES CREATED (1):\n  ShimCacheMutex\n\nPROCESSES CREATED (2):\n  C:\\WINDOWS\\system32\\HelpMe.exe\n  HelpMe.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  kernel32.dll: GetProcAddress\n\nCOMPILER AND PACKER:\n  [Packer] ASPack\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  Packers: {'PEiD': 'ASPack v2.12'}\n\nPE SECTIONS (entropy):\n  CODE         entropy=8.0\n  DATA         entropy=7.1\n  BSS          entropy=0.0\n  .idata       entropy=7.86\n  .tls         entropy=0.0\n  .rdata       entropy=0.21\n  .reloc       entropy=0.0\n  .rsrc        entropy=6.69\n  .aspack      entropy=4.93\n  .adata       entropy=0.0\n  [overlay]    entropy=5.531412124633789  size=399251\n",
  "ground_truth": {
    "hash": "3a863472603522f8123dbc92d1cfc2f606aae2ef61c66623c05856b02dda5425",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9296,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "packed with ASPack"
      },
      {
        "id": "T1547",
        "signature_description": "reference startup folder"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "contain loop",
      "contain a thread local storage (.tls) section",
      "packed with ASPack",
      "contain a resource (.rsrc) section",
      "reference startup folder",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SoftWare\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell"
      ],
      "registry_all_count": 4,
      "files_dropped": [
        "420.csv.exe",
        "664.csv.exe",
        "AUTORUN.INF.exe",
        "C:\\AUTORUN.INF.exe",
        "DUlHWC.dll.exe",
        "DefaultStore_59R.bin.exe",
        "HelpMe.exe",
        "UserMigratedStore_59R.bin.exe",
        "\\AUTORUN.INF",
        "\\AUTORUN.INF.exe",
        "\\AutoRun.exe",
        "\\MSOCache\\All Users\\{90140000-001B-0409-0000-0000000FF1CE}-C\\WordMUI.xml.exe",
        "\\Windows\\System32\\HelpMe.exe",
        "chrome.exe.exe",
        "index.dat.exe",
        "wordpfct.wpd.exe",
        "wordpfct.wpg.exe"
      ],
      "mutexes": [
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\WINDOWS\\system32\\HelpMe.exe",
        "HelpMe.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "SNDBOX",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}