# dataset/preprocessed/3b168413c0af507516877e0b51c85d954047b1455b60ffe149665a7767cdedba.json
{
  "hash": "3b168413c0af507516877e0b51c85d954047b1455b60ffe149665a7767cdedba",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  724.csv.exe\n    C:\\Windows\\system32\\HelpMe.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (2):\n  - packed with ASPack\n  - reference startup folder\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - contain loop\n  - Execution::Install Additional Program [B0023]\n  - contain a thread local storage (.tls) section\n  - Anti-Static Analysis::Software Packing [F0001]\n  - contain a resource (.rsrc) section\n  - reference startup folder\n  - (internal) packer file limitation\n  - Discovery::Analysis Tool Discovery::Process detection [B0013.001]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (13):\n  420.csv.exe  [sha256: d081356bcdb012aa...]\n  724.csv.exe  [sha256: 00ecf578a1c51a72...]\n  AUTORUN.INF.exe  [sha256: 84f6bf7913f9b705...]\n  chrome.exe.exe  [sha256: aecd1bb462fd188b...]\n  DefaultStore_59R.bin.exe  [sha256: 486027aa6daf66e0...]\n  HelpMe.exe  [sha256: 7ff3b664f0cd5a88...]\n  index.dat.exe  [sha256: aaa11b7a48633393...]\n  QcQOWN.dll.exe  [sha256: 82a2d4e50e520035...]\n  UserMigratedStore_59R.bin.exe  [sha256: 90a4617b099553d5...]\n  uyKBNp.dll.exe  [sha256: c00b3c0861285205...]\n  wordpfct.wpd.exe  [sha256: b4ebd4496d6a7c77...]\n  wordpfct.wpg.exe  [sha256: 206399c01e569315...]\n  C:\\AUTORUN.INF.exe  [sha256: 5fbde351202f56d9...]\n\nFILES WRITTEN (30):\n  C:\\AUTOEXEC.BAT.exe\n  C:\\AUTORUN.INF.exe\n  C:\\CONFIG.SYS.exe\n  C:\\Documents and Settings\\<USER>\\Start Menu\\Programs\\Startup\\Soft.lnk\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Crypto\\RSA\\S-1-5-18\\d42cc0c3858a58db2db37658219e6400_fabbc6a1-c573-4ea0-9ca1-50004b35a440.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Media Player\\DefaultStore_59R.bin.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Media Player\\UserMigratedStore_59R.bin.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Connections\\Pbk\\sharedaccess.ini.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\qmgr0.dat.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\qmgr1.dat.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\<USER>.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\Default Pictures\\airplane.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\Default Pictures\\astronaut.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\guest.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\Beethoven's Symphony No. 9 (Scherzo).wma.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\New Stories (Highway Blues).wma.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Playlists\\000511FF\\Plylst1.wpl.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SoftWare\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell = Explorer.exe  HelpMe.exe\n\nOTHER REGISTRY KEYS SET (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Explorer\\Advanced\\Folder\\Hidden\\SHOWALL\\CheckedValue\n\nMUTEXES CREATED (1):\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  C:\\WINDOWS\\system32\\HelpMe.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  kernel32.dll: GetProcAddress\n\nCOMPILER AND PACKER:\n  [Packer] ASPack\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  Packers: {'PEiD': 'ASPack v2.12'}\n\nPE SECTIONS (entropy):\n  CODE         entropy=8.0\n  DATA         entropy=7.1\n  BSS          entropy=0.0\n  .idata       entropy=7.86\n  .tls         entropy=0.0\n  .rdata       entropy=0.21\n  .reloc       entropy=0.0\n  .rsrc        entropy=6.69\n  .aspack      entropy=4.93\n  .adata       entropy=0.0\n  [overlay]    entropy=5.542478561401367  size=94273\n",
  "ground_truth": {
    "hash": "3b168413c0af507516877e0b51c85d954047b1455b60ffe149665a7767cdedba",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9167,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "packed with ASPack"
      },
      {
        "id": "T1547",
        "signature_description": "reference startup folder"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "contain loop",
      "contain an embedded PE file",
      "contain a thread local storage (.tls) section",
      "packed with ASPack",
      "contain a resource (.rsrc) section",
      "reference startup folder",
      "(internal) packer file limitation",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SoftWare\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell"
      ],
      "registry_all_count": 2,
      "files_dropped": [
        "420.csv.exe",
        "724.csv.exe",
        "AUTORUN.INF.exe",
        "C:\\AUTORUN.INF.exe",
        "DefaultStore_59R.bin.exe",
        "HelpMe.exe",
        "QcQOWN.dll.exe",
        "UserMigratedStore_59R.bin.exe",
        "chrome.exe.exe",
        "index.dat.exe",
        "uyKBNp.dll.exe",
        "wordpfct.wpd.exe",
        "wordpfct.wpg.exe"
      ],
      "mutexes": [
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\WINDOWS\\system32\\HelpMe.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}