# dataset/preprocessed/3b665d164635bb9c7c1ffadf5a9063eec20664737f8268ed27191bf485fa509c.json
{
  "hash": "3b665d164635bb9c7c1ffadf5a9063eec20664737f8268ed27191bf485fa509c",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  <PATH_SAMPLE.EXE>\n    %WINDIR%\\syswow64\\schtasks.exe\n    %WINDIR%\\syswow64\\cmd.exe\n  <SYSTEM32>\\svchost.exe\n  <PATH_SAMPLE.EXE>\n    %WINDIR%\\syswow64\\schtasks.exe\n    %WINDIR%\\syswow64\\cmd.exe\n  <PATH_SAMPLE.EXE>\n  %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\choice.exe\n  %WINDIR%\\syswow64\\cmd.exe\n    %WINDIR%\\syswow64\\choice.exe\n  %WINDIR%\\syswow64\\schtasks.exe\n  <PATH_SAMPLE.EXE>\n  %WINDIR%\\syswow64\\choice.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (32):\n  - set file attributes\n  - query or enumerate registry key\n  - query or enumerate registry value\n  - find process by name\n  - query environment variable\n  - get disk information\n  - access WMI data in .NET\n  - enumerate processes\n  - get hostname\n  - reference WMI statements\n  - get session integrity level\n  - get session user name\n  - capture screenshot\n  - get common file path\n  - delete registry key\n  - check if file exists\n  - delete registry value\n  - schedule task via schtasks\n  - encrypt data using DPAPI\n  - enumerate files on Windows\n  - gather chrome based browser login information\n  - link function at runtime on Windows\n  - check if directory exists\n  - decode data using Base64 in .NET\n  - bypass Mark of the Web\n  - log keystrokes via polling\n  - get keyboard layout\n  - find process by PID\n  - create process via WMI in .NET\n  - reference anti-VM strings targeting VirtualBox\n  - reference anti-VM strings targeting VMWare\n  - reference anti-VM strings\n\nCAPA RULE MATCHES (MBC):\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - reference absolute stream path on Windows\n  - Process::Create Process [C0017]\n  - Process::Suspend Thread [C0055]\n  - File System::Set File Attributes [C0050]\n  - get process image filename\n  - Process::Check Mutex [C0043]\n  - Process::Create Thread [C0038]\n  - Operating System::Registry::Query Registry Key [C0036.005]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - terminate process by name in .NET\n  - find process by name\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Process::Terminate Process [C0018]\n  - Discovery::System Information Discovery [E1082]\n  - enumerate drives\n  - access WMI data in .NET\n  - Operating System::Console [C0033]\n  - enumerate processes\n  - reference WMI statements\n  - Communication::HTTP Communication::Get Response [C0002.017]\n  - Command and Control::C2 Communication::Receive Data [B0030.002]\n  - unmanaged call\n  - Discovery::Application Window Discovery [E1010]\n\nHIGHLIGHTED CALLS:\n  IsDebuggerPresent\n\n--- ARTIFACTS ---\n\nOTHER REGISTRY KEYS SET (10):\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\UNCAsIntranet\n  <HKCU>\\Software\\Microsoft\\Windows\\CurrentVersion\\Internet Settings\\ZoneMap\\\\AutoDetect\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectId_\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\LruList\\CurrentLru\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectId\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\LruList\\00000000000002FD\\ObjectLru\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\ObjectTable\\E9\\_ObjectLru_\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\ObjectTable\\E9\\_FileId_\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\ObjectTable\\E9\\_Usn_\n  \\REGISTRY\\A\\{c76ecda6-fa86-11e9-a3fd-d222989d7f0b}\\DefaultObjectStore\\ObjectTable\\E9\\_UsnJournalId_\n\nMUTEXES CREATED (19):\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1228_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1276_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1748_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1868_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_200_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2176_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2272_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2316_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2356_v4.0.30319\n  \\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2360_v4.0.30319\n\nMUTEXES OPENED (17):\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{076d5de3-9426-4f39-8385-6baac4e77f0e}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{2c344653-4e84-4dc2-9e62-0d9f5532fe8a}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{40f1a774-a601-40a1-b49e-875767912112}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{42c76acd-c5b4-4bd9-ac1e-5547f1a589b3}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{52b27fc5-5ff8-4f7c-ad92-b3cced99c10a}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{8149c6ae-a1c9-4625-b629-e336501873b2}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{85e5c119-41be-4b31-bae5-4e43ca897377}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{86adc6ae-6825-4626-80ae-dddc491d7271}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{89b6ba51-bdb9-4951-821c-c3f8279f5dd9}\n  \\Sessions\\1\\BaseNamedObjects\\Global\\BFE_Notify_Event_{9cdf6815-476b-40b7-99b3-e83eb27e6bd0}\n\nPROCESSES CREATED (6):\n  %WINDIR%\\syswow64\\choice.exe\n  %WINDIR%\\syswow64\\cmd.exe\n  %WINDIR%\\syswow64\\schtasks.exe\n  <PATH_SAMPLE.EXE>\n  <SYSTEM32>\\conhost.exe\n  <SYSTEM32>\\wbem\\wmiprvse.exe\n\n--- FILE IDENTITY ---\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Library] .NET\n  [Compiler] VB.NET\n  [Compiler] VB.NET\n  [Library] .NET\n  [Linker] Microsoft Linker\n  Packers: {'PEiD': '.NET executable'}\n\nPE SECTIONS (entropy):\n  .text        entropy=4.69\n  .rsrc        entropy=3.83\n  .reloc       entropy=0.1\n",
  "ground_truth": {
    "hash": "3b665d164635bb9c7c1ffadf5a9063eec20664737f8268ed27191bf485fa509c",
    "verdict": "MALICIOUS",
    "detection_rate": 0.7746,
    "malicious_count": 55,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry key"
      },
      {
        "id": "T1027",
        "signature_description": "encrypt data using DPAPI"
      },
      {
        "id": "T1033",
        "signature_description": "get session integrity level"
      },
      {
        "id": "T1047",
        "signature_description": "access WMI data in .NET"
      },
      {
        "id": "T1053",
        "signature_description": "schedule task via schtasks"
      },
      {
        "id": "T1056",
        "signature_description": "log keystrokes via polling"
      },
      {
        "id": "T1057",
        "signature_description": "find process by name"
      },
      {
        "id": "T1082",
        "signature_description": "query environment variable"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry key"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1140",
        "signature_description": "decode data using Base64 in .NET"
      },
      {
        "id": "T1213",
        "signature_description": "reference WMI statements"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1497",
        "signature_description": "reference anti-VM strings targeting VirtualBox"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1553",
        "signature_description": "bypass Mark of the Web"
      },
      {
        "id": "T1555",
        "signature_description": "gather chrome based browser login information"
      },
      {
        "id": "T1614",
        "signature_description": "get keyboard layout"
      }
    ],
    "capa_sig_names": [
      "create or open registry key",
      "reference absolute stream path on Windows",
      "create process on Windows",
      "suspend thread",
      "set file attributes",
      "get process image filename",
      "check mutex",
      "create thread",
      "query or enumerate registry key",
      "query or enumerate registry value",
      "terminate process by name in .NET",
      "find process by name",
      "set registry value",
      "terminate process",
      "query environment variable",
      "get disk information",
      "enumerate drives",
      "access WMI data in .NET",
      "manipulate console buffer",
      "create a process with modified I/O handles and window",
      "enumerate processes",
      "get hostname",
      "reference WMI statements",
      "read data from Internet",
      "receive data",
      "unmanaged call",
      "get graphical window text",
      "get session integrity level",
      "get session user name",
      "capture screenshot",
      "save image in .NET",
      "get common file path",
      "delete registry key",
      "delete file",
      "move file",
      "check if file exists",
      "delete registry value",
      "schedule task via schtasks",
      "encrypt data using DPAPI",
      "manipulate unmanaged memory in .NET",
      "enumerate files on Windows",
      "gather chrome based browser login information",
      "set current directory",
      "find data using regex in .NET",
      "link function at runtime on Windows",
      "check if directory exists",
      "hash data using SHA1",
      "initialize hashing via WinCrypt",
      "hash data using SHA1 via WinCrypt",
      "hash data via WinCrypt",
      "load XML in .NET",
      "decode data using Base64 in .NET",
      "bypass Mark of the Web",
      "create mutex",
      "manipulate network credentials in .NET",
      "send HTTP request",
      "create HTTP request",
      "receive HTTP response",
      "send data",
      "send file using FTP",
      "contain loop",
      "log keystrokes via polling",
      "set application hook",
      "get keyboard layout",
      "find process by PID",
      "create process via WMI in .NET",
      "compiled to the .NET platform",
      "reference anti-VM strings targeting VirtualBox",
      "reference anti-VM strings targeting VMWare",
      "reference anti-VM strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [],
      "registry_all_count": 48,
      "files_dropped": [],
      "mutexes": [
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1228_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1276_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1748_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_1868_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_200_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2176_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2272_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2316_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2356_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2360_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2504_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2636_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2836_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\Global\\CPFATE_2944_v4.0.30319",
        "\\Sessions\\1\\BaseNamedObjects\\ec3a72004335486eb2e637b988301a65",
        "\\Sessions\\1\\BaseNamedObjects\\ec3a72004335486eb2e637b988301a65_ewh",
        "\\Sessions\\1\\BaseNamedObjects\\ec3a72004335486eb2e637b988301a65_logger",
        "ec3a72004335486eb2e637b988301a65",
        "ec3a72004335486eb2e637b988301a65_logger"
      ],
      "processes": [
        "%WINDIR%\\syswow64\\choice.exe",
        "%WINDIR%\\syswow64\\cmd.exe",
        "%WINDIR%\\syswow64\\schtasks.exe",
        "<PATH_SAMPLE.EXE>",
        "<SYSTEM32>\\conhost.exe",
        "<SYSTEM32>\\wbem\\wmiprvse.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "Dr.Web vxCube",
    "VirusTotal Jujubox"
  ]
}