# dataset/preprocessed/3cf264d6bdb61f66c016478eb4e37dbafbb7a34ba5e85ed9c2b7fbdb1e121f82.json
{
  "hash": "3cf264d6bdb61f66c016478eb4e37dbafbb7a34ba5e85ed9c2b7fbdb1e121f82",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  3cf264d6bdb61f66_japanese beastiality trambling public hairy .avi.exe\n    C:\\Users\\<USER>\\Downloads\\3cf264d6bdb61f66_japanese beastiality trambling public hairy .avi.exe\n    C:\\Users\\<USER>\\Downloads\\3cf264d6bdb61f66_japanese beastiality trambling public hairy .avi.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\3cf264d6bdb61f66c016478eb4e37dbafbb7a34ba5e85ed9c2b7fbdb1e121f82(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\3cf264d6bdb61f66_japanese beastiality trambling public hairy .avi.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (27):\n  american beastiality lesbian uncut young  (Jenna,Janette).avi.exe  [sha256: 287e8d5c63259292...]\n  american gang bang fucking big  (Melissa).rar.exe  [sha256: 9bfb2f77871f7586...]\n  american nude horse hot (!) titts wifey .zip.exe  [sha256: ded62be3f8291295...]\n  beast hidden cock circumcision .mpeg.exe  [sha256: 9b332be08d6c538b...]\n  black cumshot fucking hidden blondie .mpg.exe  [sha256: 80a7a258630288fd...]\n  black handjob xxx sleeping  (Jade).mpeg.exe  [sha256: fad6125b0b1b1a44...]\n  black kicking sperm voyeur 40+ .zip.exe  [sha256: 22fbf0857770e24f...]\n  black kicking xxx [bangbus] glans .zip.exe  [sha256: 96416e02b80df1f5...]\n  brasilian action hardcore several models ejaculation .mpeg.exe  [sha256: f0cd0f5d1b879962...]\n  brasilian cum sperm voyeur hole balls .rar.exe  [sha256: 31830227e76f54ac...]\n  bukkake [free] titts lady .rar.exe  [sha256: 2c2c47d9afa29bd1...]\n  danish handjob sperm licking ejaculation .avi.exe  [sha256: a11b45b097dfa829...]\n  defrag64.exe  [sha256: 1a02fd89f7d1fd5f...]\n  driver.exe  [sha256: affabe261a6051f4...]\n  drwatson32.exe  [sha256: 44cc8d0e12795cae...]\n  freeftp.exe  [sha256: 3c6de8208edd8875...]\n  fucking big sm .avi.exe  [sha256: 0f21c0123108c63d...]\n  hardcore sleeping lady .mpeg.exe  [sha256: 15924377b9e9826f...]\n  C:\\Windows\\mssrv.exe  [sha256: 8a8799a8be2cf005...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\trambling [free] legs  (Gina).rar.exe  [sha256: 93634d383b90cbed...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\russian nude bukkake girls gorgeoushorny .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\american beastiality lesbian uncut young  (Jenna,Janette).avi.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\black kicking sperm voyeur 40+ .zip.exe\n  C:\\Documents and Settings\\All Users\\Templates\\horse sleeping 50+ .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\swedish handjob xxx hidden glans blondie .avi.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\danish handjob sperm licking ejaculation .avi.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\american nude horse hot (!) titts wifey .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\russian handjob xxx hidden .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\american gang bang fucking big  (Melissa).rar.exe\n  C:\\WINDOWS\\Downloaded Program Files\\horse licking beautyfull .mpg.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\tyrkish horse lingerie uncut sweet .avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\japanese cum beast hot (!) fishy .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\spanish gay masturbation titts balls .rar.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\horse hot (!)  (Liz).avi.exe\n  C:\\WINDOWS\\assembly\\temp\\brasilian cum sperm voyeur hole balls .rar.exe\n  C:\\WINDOWS\\assembly\\tmp\\black cumshot fucking hidden blondie .mpg.exe\n  C:\\WINDOWS\\ime\\shared\\italian gang bang horse lesbian .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\brasilian action hardcore several models ejaculation .mpeg.exe\n  C:\\WINDOWS\\security\\templates\\hardcore [milf] swallow .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.879913806915283  size=395178\n",
  "ground_truth": {
    "hash": "3cf264d6bdb61f66c016478eb4e37dbafbb7a34ba5e85ed9c2b7fbdb1e121f82",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\canadian xxx licking  (Gina,Kathrin).mpg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\handjob full movie  (Gina).mpeg.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\american horse horse licking 40+  (Ashley).mpg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\trambling [free] legs  (Gina).rar.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\indian cumshot [bangbus] .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\canadian gay voyeur ejaculation .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\black cum lesbian boobs penetration  (Sonja).avi.exe",
        "C:\\Program Files\\dotnet\\shared\\russian handjob several models gorgeoushorny .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "american beastiality lesbian uncut young  (Jenna,Janette).avi.exe",
        "american gang bang fucking big  (Melissa).rar.exe",
        "american nude horse hot (!) titts wifey .zip.exe",
        "beast hidden cock circumcision .mpeg.exe",
        "black cumshot fucking hidden blondie .mpg.exe",
        "black handjob xxx sleeping  (Jade).mpeg.exe",
        "black kicking sperm voyeur 40+ .zip.exe",
        "black kicking xxx [bangbus] glans .zip.exe",
        "brasilian action hardcore several models ejaculation .mpeg.exe",
        "brasilian cum sperm voyeur hole balls .rar.exe",
        "bukkake [free] titts lady .rar.exe",
        "danish handjob sperm licking ejaculation .avi.exe",
        "defrag64.exe",
        "driver.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "fucking big sm .avi.exe",
        "hardcore sleeping lady .mpeg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}