# dataset/preprocessed/3ea179ae990cea606a776875bb3005e7de7c8a206fd942496a937f1244fbb57a.json
{
  "hash": "3ea179ae990cea606a776875bb3005e7de7c8a206fd942496a937f1244fbb57a",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  3ea179ae990cea60_british lingerie uncut nipples redhair .rar.exe\n    C:\\Users\\<USER>\\Downloads\\3ea179ae990cea60_british lingerie uncut nipples redhair .rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\3ea179ae990cea606a776875bb3005e7de7c8a206fd942496a937f1244fbb57a(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\3ea179ae990cea60_british lingerie uncut nipples redhair .rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  action hot (!) .rar.exe  [sha256: ab5e3ac411130716...]\n  action trambling voyeur .rar.exe  [sha256: c0882cb5bffb3243...]\n  asian kicking xxx uncut glans sm .avi.exe  [sha256: 9a032384d3f98ed2...]\n  beast big balls  (Curtney).mpeg.exe  [sha256: 8d5275ea1ea75057...]\n  beast [milf] .mpg.exe  [sha256: 7663d8c489026067...]\n  beastiality blowjob several models vagina black hairunshaved .avi.exe  [sha256: b3f5e308720530ba...]\n  beastiality girls nipples mistress .rar.exe  [sha256: 8602fd8fb93e9323...]\n  black horse sleeping stockings  (Sandy).avi.exe  [sha256: 7b16c5ced1d38199...]\n  british beast hardcore masturbation cock ejaculation .rar.exe  [sha256: c067f6c4312e579a...]\n  canadian fetish beastiality hot (!) hole .mpg.exe  [sha256: 249ad8c0d3e0cf82...]\n  cum hot (!) girly .zip.exe  [sha256: 415c3ce8a7b149b3...]\n  danish trambling horse masturbation .mpeg.exe  [sha256: 8c03ccd342b2990e...]\n  defrag.exe  [sha256: 999edf9d7c249daf...]\n  deinstall.exe  [sha256: 5068242602c3cf8d...]\n  edit.exe  [sha256: 9bbb356bcf65cd6a...]\n  explorer.exe  [sha256: b354df801ef4777f...]\n  freeftp.exe  [sha256: a7d7aaf07dce67ed...]\n  french action fetish sleeping hole hairy .zip.exe  [sha256: 35715721efc0839c...]\n  gang bang full movie penetration .mpg.exe  [sha256: 18cf2ad4551eca85...]\n  gay sleeping .rar.exe  [sha256: a3f034111ec1b038...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Templates\\norwegian horse licking stockings .zip.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\beastiality blowjob several models vagina black hairunshaved .avi.exe\n  C:\\Documents and Settings\\All Users\\Templates\\canadian fetish beastiality hot (!) hole .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\french action fetish sleeping hole hairy .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\beastiality girls nipples mistress .rar.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\action trambling voyeur .rar.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\indian cum big ejaculation .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\british beast hardcore masturbation cock ejaculation .rar.exe\n  C:\\Program Files\\Movie Maker\\Shared\\danish trambling horse masturbation .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\gay sleeping .rar.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\sperm catfight girly .avi.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\swedish trambling [milf] ash bedroom .mpg.exe\n  C:\\WINDOWS\\Temp\\action hot (!) .rar.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\black horse sleeping stockings  (Sandy).avi.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\porn girls beautyfull  (Sonja).mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\indian fetish [milf] .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\kicking trambling several models titts castration .avi.exe\n  C:\\WINDOWS\\ime\\shared\\german cum licking hole high heels .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\japanese beast nude voyeur  (Jenna).mpeg.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.87\n  [overlay]    entropy=7.879909038543701  size=1248781\n",
  "ground_truth": {
    "hash": "3ea179ae990cea606a776875bb3005e7de7c8a206fd942496a937f1244fbb57a",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9577,
    "malicious_count": 68,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\gay kicking hot (!) vagina girly  (Sylvia,Anniston).zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\malaysia bukkake trambling masturbation hole .mpg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\danish trambling licking beautyfull .rar.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\fucking beastiality lesbian 40+ .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\british beast horse [milf] cock circumcision .zip.exe",
        "C:\\Program Files\\dotnet\\shared\\danish beast nude sleeping redhair  (Melissa).rar.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "action hot (!) .rar.exe",
        "action trambling voyeur .rar.exe",
        "asian kicking xxx uncut glans sm .avi.exe",
        "beast [milf] .mpg.exe",
        "beast big balls  (Curtney).mpeg.exe",
        "beastiality blowjob several models vagina black hairunshaved .avi.exe",
        "beastiality girls nipples mistress .rar.exe",
        "black horse sleeping stockings  (Sandy).avi.exe",
        "british beast hardcore masturbation cock ejaculation .rar.exe",
        "canadian fetish beastiality hot (!) hole .mpg.exe",
        "cum hot (!) girly .zip.exe",
        "danish trambling horse masturbation .mpeg.exe",
        "defrag.exe",
        "deinstall.exe",
        "edit.exe",
        "explorer.exe",
        "freeftp.exe",
        "french action fetish sleeping hole hairy .zip.exe",
        "gang bang full movie penetration .mpg.exe",
        "gay sleeping .rar.exe",
        "german cum licking hole high heels .zip.exe",
        "indian beast handjob [free] castration  (Melissa,Sarah).avi.exe",
        "indian cum big ejaculation .avi.exe",
        "indian fetish [milf] .mpeg.exe",
        "install.exe",
        "japanese beast nude voyeur  (Jenna).mpeg.exe",
        "kicking trambling several models titts castration .avi.exe",
        "mssrv.exe",
        "norwegian horse licking stockings .zip.exe",
        "notes.txt.exe",
        "porn girls beautyfull  (Sonja).mpg.exe",
        "rundll64.exe",
        "russian bukkake lesbian [bangbus] ash .zip.exe",
        "safetyserver.exe",
        "sperm catfight girly .avi.exe",
        "swedish animal kicking licking cock latex  (Karin,Kathrin).rar.exe",
        "swedish trambling [milf] ash bedroom .mpg.exe",
        "tyrkish horse gang bang several models  (Curtney,Janette).avi.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}