# dataset/preprocessed/3f22b72b633f626f9e501a9a29e78c139be5a1280254b086f4120d11fb28eaa9.json
{
  "hash": "3f22b72b633f626f9e501a9a29e78c139be5a1280254b086f4120d11fb28eaa9",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  3f22b72b633f626f_mssrv.exe\n    C:\\Users\\<USER>\\Downloads\\3f22b72b633f626f_mssrv.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\3f22b72b633f626f9e501a9a29e78c139be5a1280254b086f4120d11fb28eaa9(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\3f22b72b633f626f_mssrv.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (19):\n  action catfight boobs black hairunshaved .mpeg.exe  [sha256: 01759ead3c5789c8...]\n  african cum hidden .zip.exe  [sha256: c6a489b7673e7efa...]\n  african nude catfight .avi.exe  [sha256: 8ad6efe1960aca55...]\n  asian beast catfight swallow  (Janette).mpg.exe  [sha256: af5b441b0c15ce38...]\n  black beastiality hardcore public beautyfull  (Kathrin,Kathrin).mpg.exe  [sha256: 055aa643c6a7806d...]\n  blowjob big cock  (Sonja).rar.exe  [sha256: 3280a515b9df08ba...]\n  brasilian fetish catfight lady .zip.exe  [sha256: bee74e176e49109f...]\n  brasilian fetish cumshot full movie .avi.exe  [sha256: e70d84cdfea69157...]\n  defrag64.exe  [sha256: 108110d6019d3cf7...]\n  deinstall.exe  [sha256: 262dd18de08bd2a0...]\n  driver.exe  [sha256: 2fad3efef7085921...]\n  drwatson32.exe  [sha256: cb8c4fc6881b46f5...]\n  C:\\Windows\\mssrv.exe  [sha256: bc1592c367c991fc...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\kicking public titts high heels  (Curtney).zip.exe  [sha256: 8dd295960c3fce68...]\n  C:\\Program Files\\dotnet\\shared\\norwegian bukkake lesbian boobs  (Samantha).rar.exe  [sha256: 94c68ee50c5ea621...]\n  C:\\Program Files\\Windows Journal\\Templates\\french fucking fetish several models .mpg.exe  [sha256: adfd3047d5621034...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\japanese trambling uncut vagina mistress .rar.exe  [sha256: f116742c279fc0a2...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\american fucking handjob [free] .mpeg.exe  [sha256: c2338cf431ff2728...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\norwegian horse lingerie hidden young  (Christine).avi.exe  [sha256: 9e92c7df3de08a0a...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\horse [free] ash sweet .mpeg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\fetish masturbation .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\black beastiality hardcore public beautyfull  (Kathrin,Kathrin).mpg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\brasilian fetish cumshot full movie .avi.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\lesbian voyeur mature .avi.exe\n  C:\\Documents and Settings\\Default User\\Templates\\action catfight boobs black hairunshaved .mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\trambling bukkake hidden high heels  (Christine,Kathrin).rar.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\french fucking lingerie [milf] .mpeg.exe\n  C:\\Program Files\\Movie Maker\\Shared\\russian kicking sleeping .mpeg.exe\n  C:\\WINDOWS\\Downloaded Program Files\\italian cum [milf] 40+ .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\fetish hot (!) nipples .rar.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\japanese porn sperm big girly  (Jenna,Gina).avi.exe\n  C:\\WINDOWS\\Temp\\malaysia animal lesbian girls ash .mpeg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\gang bang sperm public shower .mpeg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\sperm [milf] .mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\spanish fetish [bangbus] .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\hardcore voyeur swallow .zip.exe\n  C:\\WINDOWS\\ime\\shared\\spanish horse hidden .mpeg.exe\n  C:\\WINDOWS\\mssrv.exe\n  C:\\WINDOWS\\pchealth\\helpctr\\Temp\\spanish nude big leather .rar.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=4.86\n  [overlay]    entropy=7.879809379577637  size=421686\n",
  "ground_truth": {
    "hash": "3f22b72b633f626f9e501a9a29e78c139be5a1280254b086f4120d11fb28eaa9",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9437,
    "malicious_count": 67,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\american fucking handjob [free] .mpeg.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\norwegian horse lingerie hidden young  (Christine).avi.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\kicking public titts high heels  (Curtney).zip.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\french fucking fetish several models .mpg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\japanese trambling uncut vagina mistress .rar.exe",
        "C:\\Program Files\\dotnet\\shared\\norwegian bukkake lesbian boobs  (Samantha).rar.exe",
        "C:\\Windows\\mssrv.exe",
        "action catfight boobs black hairunshaved .mpeg.exe",
        "african cum hidden .zip.exe",
        "african nude catfight .avi.exe",
        "asian beast catfight swallow  (Janette).mpg.exe",
        "black beastiality hardcore public beautyfull  (Kathrin,Kathrin).mpg.exe",
        "blowjob big cock  (Sonja).rar.exe",
        "brasilian fetish catfight lady .zip.exe",
        "brasilian fetish cumshot full movie .avi.exe",
        "defrag64.exe",
        "deinstall.exe",
        "driver.exe",
        "drwatson32.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}