# dataset/preprocessed/3f274dfe5f938c60da8b53fe106063105a9bc0bd0960280314f74fba5e466d80.json
{
  "hash": "3f274dfe5f938c60da8b53fe106063105a9bc0bd0960280314f74fba5e466d80",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  3f274dfe5f938c60_lingerie voyeur hole traffic  (janette).rar.exe\n    C:\\Users\\<USER>\\Downloads\\3f274dfe5f938c60_lingerie voyeur hole traffic  (janette).rar.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\3f274dfe5f938c60da8b53fe106063105a9bc0bd0960280314f74fba5e466d80(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\3f274dfe5f938c60_lingerie voyeur hole traffic  (janette).rar.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american porn xxx big titts circumcision .avi.exe  [sha256: 7af09192bcc148af...]\n  black cumshot lesbian hidden feet .zip.exe  [sha256: 9c9757a5c0df57af...]\n  black cumshot sperm uncut feet .zip.exe  [sha256: 3ed934d45e809f6f...]\n  black nude blowjob voyeur latex  (Kathrin,Curtney).mpeg.exe  [sha256: 8f5d2bdd938ee1ce...]\n  blowjob masturbation cock .avi.exe  [sha256: 162200422927aa26...]\n  brasilian beastiality lingerie several models glans fishy .avi.exe  [sha256: 55e9e52bd4d4220c...]\n  brasilian cumshot xxx lesbian feet .mpg.exe  [sha256: fa5c3aef0b808837...]\n  brasilian fetish blowjob girls cock fishy .mpeg.exe  [sha256: 7592431ea68f626f...]\n  brasilian nude fucking big glans .zip.exe  [sha256: 2f42225f0b9921bd...]\n  brasilian porn sperm uncut redhair .mpeg.exe  [sha256: 2e6f6011175fbb85...]\n  bukkake [bangbus] 40+ .mpeg.exe  [sha256: 4e22056b01e8948e...]\n  defrag.exe  [sha256: e998ac1b8c7bca8f...]\n  deinstall.exe  [sha256: 6b7fbad8840916e7...]\n  drwatson32.exe  [sha256: 9b26c5823ada06fe...]\n  freeftp.exe  [sha256: e1b21f381957e3d2...]\n  indian animal sperm hidden fishy  (Kathrin,Samantha).mpg.exe  [sha256: efe6f3e2a19cc56e...]\n  indian gang bang blowjob catfight granny  (Jenna,Melissa).avi.exe  [sha256: 85fa96fe016ba19b...]\n  indian nude horse lesbian cock circumcision .mpeg.exe  [sha256: f2e1f545d9c2ff3c...]\n  italian beastiality blowjob sleeping glans black hairunshaved .zip.exe  [sha256: 475f6ff502c08d4f...]\n  lesbian hot (!) feet 40+  (Curtney).mpeg.exe  [sha256: 6423a221c99d21a7...]\n\nFILES WRITTEN (30):\n  C:\\3f274dfe5f938c60da8b53fe106063105a9bc0bd0960280314f74fba5e466d80\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\tyrkish gang bang hardcore catfight .avi.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\black nude blowjob voyeur latex  (Kathrin,Curtney).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\bukkake [bangbus] 40+ .mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\tyrkish horse lingerie public boots .mpg.exe\n  C:\\Documents and Settings\\Default User\\Templates\\russian nude sperm uncut lady .zip.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\black cumshot lesbian hidden feet .zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\sperm masturbation  (Jade).zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\russian animal trambling hot (!) black hairunshaved .avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\trambling [free] beautyfull .zip.exe\n  C:\\WINDOWS\\Downloaded Program Files\\italian beastiality blowjob sleeping glans black hairunshaved .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\brasilian nude fucking big glans .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\tyrkish animal trambling girls titts .rar.exe\n  C:\\WINDOWS\\Temp\\american porn xxx big titts circumcision .avi.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\tyrkish horse bukkake [bangbus] sm .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\indian gang bang blowjob catfight granny  (Jenna,Melissa).avi.exe\n  C:\\WINDOWS\\assembly\\temp\\indian nude horse lesbian cock circumcision .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\indian animal sperm hidden fishy  (Kathrin,Samantha).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\trambling [milf] latex .zip.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880120277404785  size=1932305\n",
  "ground_truth": {
    "hash": "3f274dfe5f938c60da8b53fe106063105a9bc0bd0960280314f74fba5e466d80",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9286,
    "malicious_count": 65,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\japanese handjob xxx public cock  (Sonja,Jade).mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\black action horse lesbian hotel .mpeg.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\american handjob hardcore masturbation blondie .mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\russian gang bang hardcore girls swallow .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\italian handjob gay uncut 50+ .avi.exe",
        "C:\\Program Files\\dotnet\\shared\\japanese handjob horse sleeping .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "american porn xxx big titts circumcision .avi.exe",
        "black cumshot lesbian hidden feet .zip.exe",
        "black cumshot sperm uncut feet .zip.exe",
        "black nude blowjob voyeur latex  (Kathrin,Curtney).mpeg.exe",
        "blowjob masturbation cock .avi.exe",
        "brasilian beastiality lingerie several models glans fishy .avi.exe",
        "brasilian cumshot xxx lesbian feet .mpg.exe",
        "brasilian fetish blowjob girls cock fishy .mpeg.exe",
        "brasilian nude fucking big glans .zip.exe",
        "brasilian porn sperm uncut redhair .mpeg.exe",
        "bukkake [bangbus] 40+ .mpeg.exe",
        "defrag.exe",
        "deinstall.exe",
        "drwatson32.exe",
        "freeftp.exe",
        "indian animal sperm hidden fishy  (Kathrin,Samantha).mpg.exe",
        "indian gang bang blowjob catfight granny  (Jenna,Melissa).avi.exe",
        "indian nude horse lesbian cock circumcision .mpeg.exe",
        "italian beastiality blowjob sleeping glans black hairunshaved .zip.exe",
        "lesbian hot (!) feet 40+  (Curtney).mpeg.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "russian animal trambling hot (!) black hairunshaved .avi.exe",
        "russian nude sperm uncut lady .zip.exe",
        "screensaver.scr",
        "setup.exe",
        "soccer.exe",
        "sperm masturbation  (Jade).zip.exe",
        "swedish handjob xxx uncut sweet  (Sonja,Tatjana).mpeg.exe",
        "trambling [free] beautyfull .zip.exe",
        "trambling [milf] latex .zip.exe",
        "tyrkish animal trambling girls titts .rar.exe",
        "tyrkish cumshot lingerie several models cock stockings  (Samantha).mpeg.exe",
        "tyrkish gang bang hardcore catfight .avi.exe",
        "tyrkish horse bukkake [bangbus] sm .mpg.exe",
        "tyrkish horse lingerie public boots .mpg.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}