# dataset/preprocessed/3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a.json
{
  "hash": "3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  718452dc8ef422d1e4235ec6dba963b7.exe\n    c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202.exe\n      c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202a.exe\n        c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202b.exe\n          c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202c.exe\n            c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202d.exe\n              c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202e.exe\n                c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202f.exe\n                  c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202g.exe\n                    c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202h.exe\n                      c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202i.exe\n                        c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202j.exe\n                          c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202k.exe\n                            c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202l.exe\n                              c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202m.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (25):\n  - encode data using XOR\n  - hide graphical window\n  - set file attributes\n  - access PEB ldr_data\n  - get disk information\n  - get common file path\n  - persist via Run registry key\n  - find graphical window\n  - enumerate files on Windows\n  - query environment variable\n  - query or enumerate registry value\n  - delete registry value\n  - query or enumerate registry key\n  - delete registry key\n  - get memory capacity\n  - enumerate files recursively\n  - get file size\n  - check if file exists\n  - check OS version\n  - get disk size\n  - link function at runtime on Windows\n  - get file version info\n  - link many functions at runtime\n  - enumerate gui resources\n  - packed with UPX\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - hide graphical window\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - File System::Get File Attributes [C0049]\n  - File System::Set File Attributes [C0050]\n  - Process::Create Process [C0017]\n  - Anti-Behavioral Analysis::Debugger Detection::Process Environment Block [B0001.019]\n  - Process::Check Mutex [C0043]\n  - reference absolute stream path on Windows\n  - Process::Create Thread [C0038]\n  - access PEB ldr_data\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Discovery::File and Directory Discovery [E1083]\n  - resolve function by parsing PE exports\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - find graphical window\n  - File System::Writes File [C0052]\n  - File System::Copy File [C0045]\n  - File System::Read File [C0051]\n\n--- ARTIFACTS ---\n\nFILES DROPPED (1):\n  3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a_3202.exe  [sha256: 9d628a940d31e55f...]\n\nFILES WRITTEN (1):\n  c:\\3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Trickler = \"c:\\users\\<USER>\\downloads\\718452dc8ef422d1e4235ec6dba963b7_3202.exe\"\n\nOTHER REGISTRY KEYS SET (3):\n  HKLM\\Software\\CLASSES\\CLSID\\{21FFB6C0-0DA1-11D5-A9D5-00500413153C}\\uets\n  HKLM\\software\\Qwertyuio\\Trickler\\AppPath\n  HKLM\\software\\Qwertyuio\\Trickler\\OldTrickler\n\nMUTEXES CREATED (1):\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  c:\\3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a c:\\3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a _3202.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.DLL: CreateProcessA, GetProcAddress, VirtualAlloc\n\nCERTIFICATE INFO:\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++'}\n\nPE SECTIONS (entropy):\n  UPX0         entropy=6.46\n  UPX1         entropy=4.02\n  .rsrc        entropy=3.36\n  .imports     entropy=4.54\n  [overlay]    entropy=1.7211765050888062  size=204833\n",
  "ground_truth": {
    "hash": "3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a",
    "verdict": "MALICIOUS",
    "detection_rate": 0.8611,
    "malicious_count": 62,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "gator",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "find graphical window"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1112",
        "signature_description": "delete registry value"
      },
      {
        "id": "T1129",
        "signature_description": "access PEB ldr_data"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "encode data using XOR",
      "hide graphical window",
      "delay execution",
      "create or open file",
      "create or open registry key",
      "get file attributes",
      "set file attributes",
      "create process on Windows",
      "PEB access",
      "check mutex",
      "reference absolute stream path on Windows",
      "create thread",
      "access PEB ldr_data",
      "contain loop",
      "get disk information",
      "get common file path",
      "resolve function by parsing PE exports",
      "persist via Run registry key",
      "find graphical window",
      "write file on Windows",
      "copy file",
      "read .ini file",
      "enumerate files on Windows",
      "get graphical window text",
      "query environment variable",
      "query or enumerate registry value",
      "set registry value",
      "delete registry value",
      "query or enumerate registry key",
      "delete registry key",
      "get memory capacity",
      "delete file",
      "delete directory",
      "enumerate files recursively",
      "create directory",
      "get file size",
      "check if file exists",
      "get OS version",
      "check OS version",
      "move file",
      "read file on Windows",
      "get disk size",
      "link function at runtime on Windows",
      "create mutex",
      "check mutex and exit",
      "check for time delay via GetTickCount",
      "get file version info",
      "link many functions at runtime",
      "check for software breakpoints",
      "send data on socket",
      "send data",
      "receive data on socket",
      "receive data",
      "print debug messages",
      "enumerate gui resources",
      "packed with UPX",
      "contain a resource (.rsrc) section",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\Trickler"
      ],
      "registry_all_count": 4,
      "files_dropped": [
        "3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a_3202.exe"
      ],
      "mutexes": [
        "ShimCacheMutex"
      ],
      "processes": [
        "c:\\3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a c:\\3fd27a7b678378f01596e0251012769a35605f9682e44a03cf682dc892e0891a _3202.exe"
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}