# dataset/preprocessed/4044aa2b5e79b1bc1e8997947d9ecf054fb70f56610928685d8ef24c5e757534.json
{
  "hash": "4044aa2b5e79b1bc1e8997947d9ecf054fb70f56610928685d8ef24c5e757534",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  C:\\AUTORUN.INF.exe\n    C:\\Windows\\system32\\HelpMe.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (2):\n  - packed with ASPack\n  - reference startup folder\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - contain loop\n  - contain a resource (.rsrc) section\n  - Anti-Static Analysis::Software Packing [F0001]\n  - contain a thread local storage (.tls) section\n  - Execution::Install Additional Program [B0023]\n  - reference startup folder\n  - (internal) packer file limitation\n\n--- ARTIFACTS ---\n\nFILES DROPPED (12):\n  1744.csv.exe  [sha256: 98f6fe25e27808dc...]\n  1908.csv.exe  [sha256: 25670e512771ea55...]\n  AUTORUN.INF.exe  [sha256: 79a5cc53b9aafed6...]\n  chrome.exe.exe  [sha256: e7a04433fb928e14...]\n  DefaultStore_59R.bin.exe  [sha256: 4ccc35173148a991...]\n  HelpMe.exe  [sha256: 51a5a9eed32650df...]\n  index.dat.exe  [sha256: e8e1bd70bde94806...]\n  nugidv.dll.exe  [sha256: f3a5473ab1f68951...]\n  UserMigratedStore_59R.bin.exe  [sha256: 80b363c6f6026132...]\n  wordpfct.wpd.exe  [sha256: cba10eb7842fa717...]\n  wordpfct.wpg.exe  [sha256: d7f26754785f8532...]\n  WXaUre.dll.exe  [sha256: 779df538ca7125fe...]\n\nFILES WRITTEN (30):\n  C:\\AUTOEXEC.BAT.exe\n  C:\\AUTORUN.INF.exe\n  C:\\CONFIG.SYS.exe\n  C:\\Documents and Settings\\<USER>\\Start Menu\\Programs\\Startup\\Soft.lnk\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Crypto\\RSA\\S-1-5-18\\d42cc0c3858a58db2db37658219e6400_fabbc6a1-c573-4ea0-9ca1-50004b35a440.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Media Player\\DefaultStore_59R.bin.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Media Player\\UserMigratedStore_59R.bin.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Connections\\Pbk\\sharedaccess.ini.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\qmgr0.dat.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\qmgr1.dat.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\<USER>.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\Default Pictures\\airplane.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\Default Pictures\\astronaut.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\User Account Pictures\\guest.bmp.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\Beethoven's Symphony No. 9 (Scherzo).wma.exe\n  C:\\Documents and Settings\\All Users\\Documents\\My Music\\Sample Music\\desktop.ini.exe\n  C:\\Documents and Settings\\All Users\\Documents\\desktop.ini.exe\n  C:\\WINDOWS\\system32\\HelpMe.exe\n\nPERSISTENCE REGISTRY KEYS (2):\n  HKLM\\SoftWare\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell = Explorer.exe  HelpMe.exe\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell = Explorer.exe  HelpMe.exe\u0000\n\nOTHER REGISTRY KEYS SET (2):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Explorer\\Advanced\\Folder\\Hidden\\SHOWALL\\CheckedValue\n  \\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Explorer\\Advanced\\Folder\\Hidden\\SHOWALL\\CheckedValue\n\nMUTEXES CREATED (10):\n  CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500\n  MSCTF.Shared.MUTEX.EBH\n  MSCTF.Shared.MUTEX.IMF\n  MSCTF.Shared.MUTEX.MPF\n  ShimCacheMutex\n\nPROCESSES CREATED (1):\n  C:\\WINDOWS\\system32\\HelpMe.exe\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  kernel32.dll: GetProcAddress\n\nCOMPILER AND PACKER:\n  [Packer] ASPack\n  [Compiler] Borland Delphi\n  [Linker] Turbo Linker\n  Packers: {'PEiD': 'ASPack v2.12'}\n\nPE SECTIONS (entropy):\n  CODE         entropy=8.0\n  DATA         entropy=7.1\n  BSS          entropy=0.0\n  .idata       entropy=7.86\n  .tls         entropy=0.0\n  .rdata       entropy=0.21\n  .reloc       entropy=0.0\n  .rsrc        entropy=6.69\n  .aspack      entropy=4.93\n  .adata       entropy=0.0\n  [overlay]    entropy=5.444758892059326  size=299488\n",
  "ground_truth": {
    "hash": "4044aa2b5e79b1bc1e8997947d9ecf054fb70f56610928685d8ef24c5e757534",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9296,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1027",
        "signature_description": "packed with ASPack"
      },
      {
        "id": "T1547",
        "signature_description": "reference startup folder"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "contain loop",
      "contain a resource (.rsrc) section",
      "packed with ASPack",
      "contain a thread local storage (.tls) section",
      "contain an embedded PE file",
      "reference startup folder",
      "(internal) packer file limitation"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SoftWare\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell",
        "\\REGISTRY\\MACHINE\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell"
      ],
      "registry_all_count": 4,
      "files_dropped": [
        "1744.csv.exe",
        "1908.csv.exe",
        "AUTORUN.INF.exe",
        "DefaultStore_59R.bin.exe",
        "HelpMe.exe",
        "UserMigratedStore_59R.bin.exe",
        "WXaUre.dll.exe",
        "chrome.exe.exe",
        "index.dat.exe",
        "nugidv.dll.exe",
        "wordpfct.wpd.exe",
        "wordpfct.wpg.exe"
      ],
      "mutexes": [
        "CTF.Asm.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Compart.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.LBES.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.Layouts.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TMD.MutexDefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "CTF.TimListCache.FMPDefaultS-1-5-21-1482476501-1645522239-1417001333-500MUTEX.DefaultS-1-5-21-1482476501-1645522239-1417001333-500",
        "MSCTF.Shared.MUTEX.EBH",
        "MSCTF.Shared.MUTEX.IMF",
        "MSCTF.Shared.MUTEX.MPF",
        "ShimCacheMutex"
      ],
      "processes": [
        "C:\\WINDOWS\\system32\\HelpMe.exe"
      ]
    }
  },
  "num_sandboxes": 4,
  "sandbox_names": [
    "CAPA",
    "Tencent HABO",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}