# dataset/preprocessed/4156d2efb1b232963508a472918c40c4def1ea24025884333846c50c3abce5bc.json
{
  "hash": "4156d2efb1b232963508a472918c40c4def1ea24025884333846c50c3abce5bc",
  "report_text": "--- EXECUTION CHAIN ---\n\nPROCESS TREE (flat):\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe\" \"-\"\"\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe-\"\" [C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\]\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\edhts.exe \"-\"\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (27):\n  - set file attributes\n  - encode data using XOR\n  - hide graphical window\n  - get socket status\n  - enumerate gui resources\n  - enumerate files on Windows\n  - enumerate files recursively\n  - get disk information\n  - get session user name\n  - get common file path\n  - get system information on Windows\n  - link function at runtime on Windows\n  - parse PE header\n  - link many functions at runtime\n  - get hostname\n  - capture screenshot\n  - get file size\n  - modify service\n  - check if file exists\n  - timestomp file\n  - encrypt data using RC4 PRGA\n  - get file version info\n  - check for unmoving mouse cursor\n  - accept command line arguments\n  - get geographical location\n  - check for VM using instruction VPCEXT\n  - reference the VMWare IO port\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - File System::Create File [C0016]\n  - File System::Set File Attributes [C0050]\n  - Process::Create Thread [C0038]\n  - Process::Create Process [C0017]\n  - Process::Check Mutex [C0043]\n  - Defense Evasion::Obfuscated Files or Information::Encoding-Standard Algorithm [E1027.m02]\n  - Data::Encode Data::XOR [C0026.002]\n  - File System::Get File Attributes [C0049]\n  - hide graphical window\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - File System::Copy File [C0045]\n  - File System::Delete File [C0047]\n  - Communication::Socket Communication::Send Data [C0001.007]\n  - Command and Control::C2 Communication::Send Data [B0030.001]\n  - Communication::Socket Communication::Receive Data [C0001.006]\n  - Communication::Socket Communication::Get Socket Status [C0001.012]\n  - Command and Control::C2 Communication::Receive Data [B0030.002]\n  - Anti-Behavioral Analysis::Debugger Detection::Timing/Delay Check GetTickCount [B0001.032]\n  - Process::Terminate Process [C0018]\n  - enumerate gui resources\n\n--- ARTIFACTS ---\n\nFILES DROPPED (1):\n  cagrt.exe  [sha256: 2569ea19734f486d...]\n\nFILES WRITTEN (2):\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe\n  C:\\Users\\<USER>\\AppData\\Local\\Temp\\edhts.exe\n\nPERSISTENCE REGISTRY KEYS (9):\n  HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\tzkdjynqlpts = itjhsmgooxgksuydw.exe .\n  HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\wbldiwkmgjm = itjhsmgooxgksuydw.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell = Explorer.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\\bfofjwjkdf = cphhuqmwyjuakoubwna.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\\cdjxyis = C:\\Users\\<USER>\\AppData\\Local\\Temp\\rdutfavefpzenqvbvl.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\iltjmykkc = etnpecamqdqykqyhexmgi.exe .\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\szlfmcswsxccg = C:\\Users\\<USER>\\AppData\\Local\\Temp\\rdutfavefpzenqvbvl.exe .\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\rtaprcnm = blaxhatazhpszadh.exe\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\tbojrizebhnots = C:\\Users\\<USER>\\AppData\\Local\\Temp\\rdutfavefpzenqvbvl.exe\n\nOTHER REGISTRY KEYS SET (10):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableLUA\n  HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\DisableRegistryTools\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\ConsentPromptBehaviorAdmin\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\ConsentPromptBehaviorUser\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableInstallerDetection\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableSecureUIAPaths\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\EnableVirtualization\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\PromptOnSecureDesktop\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\ValidateAdminCodeSignatures\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\System\\FilterAdministratorToken\n\nMUTEXES CREATED (5):\n  cqmnfeanpdaxsagdygagy\n  mcadxywlpfedaksroyucwfy\n  ofbfwwwkqfuesakvupgcgxx\n  tifhaaxlodbzveljfojqjr\n  zbbozbbonyvhdrgzbbozbbonyv\n\nMUTEXES OPENED (1):\n  etnpecamqdqykqyhexmgi\n\nPROCESSES CREATED (1):\n  C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe\" \"-\"\"\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: CreateProcessA, GetProcAddress, VirtualAlloc, VirtualProtect\n\nCOMPILER AND PACKER:\n  [Compiler] EP:Microsoft Visual C/C++\n  [Compiler] Microsoft Visual C/C++\n  [Linker] Microsoft Linker\n  [Tool] Visual Studio\n  Packers: {'PEiD': 'Microsoft Visual C++ v7.0'}\n\nPE SECTIONS (entropy):\n  .text        entropy=6.51\n  .rdata       entropy=5.05\n  .data        entropy=7.39\n",
  "ground_truth": {
    "hash": "4156d2efb1b232963508a472918c40c4def1ea24025884333846c50c3abce5bc",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9412,
    "malicious_count": 64,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1016",
        "signature_description": "get socket status"
      },
      {
        "id": "T1027",
        "signature_description": "encode data using XOR"
      },
      {
        "id": "T1033",
        "signature_description": "get session user name"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1070",
        "signature_description": "timestomp file"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "enumerate files on Windows"
      },
      {
        "id": "T1087",
        "signature_description": "get session user name"
      },
      {
        "id": "T1113",
        "signature_description": "capture screenshot"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1222",
        "signature_description": "set file attributes"
      },
      {
        "id": "T1497",
        "signature_description": "check for unmoving mouse cursor"
      },
      {
        "id": "T1543",
        "signature_description": "modify service"
      },
      {
        "id": "T1564",
        "signature_description": "hide graphical window"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      },
      {
        "id": "T1614",
        "signature_description": "get geographical location"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "delay execution",
      "create or open file",
      "set file attributes",
      "create thread",
      "create process on Windows",
      "check mutex",
      "encode data using XOR",
      "get file attributes",
      "hide graphical window",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "read file on Windows",
      "write file on Windows",
      "copy file",
      "delete file",
      "send data on socket",
      "send data",
      "receive data on socket",
      "get socket status",
      "receive data",
      "check for time delay via GetTickCount",
      "terminate process",
      "enumerate gui resources",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get disk information",
      "get session user name",
      "get common file path",
      "get system information on Windows",
      "get OS version",
      "link function at runtime on Windows",
      "extract resource via kernel32 functions",
      "parse PE header",
      "link many functions at runtime",
      "get hostname",
      "create mutex",
      "capture screenshot",
      "get file size",
      "write clipboard data",
      "replace clipboard data",
      "modify service",
      "check if file exists",
      "timestomp file",
      "move file",
      "create directory",
      "get graphical window text",
      "encrypt data using RC4 PRGA",
      "get file version info",
      "check for unmoving mouse cursor",
      "accept command line arguments",
      "get geographical location",
      "set environment variable",
      "create or open file",
      "read file on Windows",
      "write file on Windows",
      "copy file",
      "get common file path",
      "get OS version",
      "extract resource via kernel32 functions",
      "check for VM using instruction VPCEXT",
      "reference the VMWare IO port",
      "write clipboard data"
    ],
    "iocs": {
      "ips": [],
      "domains": [],
      "urls": [],
      "registry_persistence": [
        "HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\tzkdjynqlpts",
        "HKCU\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\wbldiwkmgjm",
        "HKLM\\SOFTWARE\\Microsoft\\Windows NT\\CurrentVersion\\Winlogon\\Shell",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\\bfofjwjkdf",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Policies\\Explorer\\Run\\cdjxyis",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\iltjmykkc",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\RunOnce\\szlfmcswsxccg",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\rtaprcnm",
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\tbojrizebhnots"
      ],
      "registry_all_count": 27,
      "files_dropped": [
        "cagrt.exe"
      ],
      "mutexes": [
        "cqmnfeanpdaxsagdygagy",
        "mcadxywlpfedaksroyucwfy",
        "ofbfwwwkqfuesakvupgcgxx",
        "tifhaaxlodbzveljfojqjr",
        "zbbozbbonyvhdrgzbbozbbonyv"
      ],
      "processes": [
        "C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe C:\\DOCUME~1\\<USER>~1\\LOCALS~1\\Temp\\cagrt.exe\" \"-\"\""
      ]
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}