# dataset/preprocessed/41fd75ae1ced9aa2f43fb88c97c1da0819b1a10ac9f16fe1e5df0f7efee3e792.json
{
  "hash": "41fd75ae1ced9aa2f43fb88c97c1da0819b1a10ac9f16fe1e5df0f7efee3e792",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (1):\n  VBOXSVR.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  41fd75ae1ced9aa2_fucking lesbian gorgeoushorny .mpg.exe\n    C:\\Users\\<USER>\\Downloads\\41fd75ae1ced9aa2_fucking lesbian gorgeoushorny .mpg.exe\n    C:\\Users\\<USER>\\Downloads\\41fd75ae1ced9aa2_fucking lesbian gorgeoushorny .mpg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\41fd75ae1ced9aa2f43fb88c97c1da0819b1a10ac9f16fe1e5df0f7efee3e792(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\41fd75ae1ced9aa2_fucking lesbian gorgeoushorny .mpg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - contain loop\n  - Discovery::System Information Discovery [E1082]\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (21):\n  action sleeping black hairunshaved .mpeg.exe  [sha256: 5afd5147739ac0ab...]\n  african porn hidden hairy .mpg.exe  [sha256: 149cfdf2fc6ccb9e...]\n  american lesbian porn hidden latex .mpeg.exe  [sha256: db8c01479476b42c...]\n  black bukkake big boobs .zip.exe  [sha256: 660dcfbf8ecbb05e...]\n  black horse [bangbus] young  (Melissa,Gina).mpeg.exe  [sha256: 194c8895fd6da8f5...]\n  brasilian bukkake big stockings  (Jenna).avi.exe  [sha256: 758a0b4ba3b5ae08...]\n  C:\\Windows\\mssrv.exe  [sha256: e4178781507ef530...]\n  C:\\Program Files\\Common Files\\Microsoft Shared\\xxx licking wifey .mpg.exe  [sha256: 135ae1d4cabcb8b8...]\n  C:\\Program Files\\dotnet\\shared\\italian lingerie [milf] mistress .avi.exe  [sha256: 8d8417c8238592f1...]\n  C:\\Program Files\\Windows Journal\\Templates\\french nude uncut legs hairy  (Ashley).mpeg.exe  [sha256: b8c7f1a8b44285b8...]\n  C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\asian handjob licking gorgeoushorny .mpeg.exe  [sha256: 2c5ed531c284840a...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\italian hardcore big .avi.exe  [sha256: 5febe8243b17216f...]\n  C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\canadian beast porn public vagina shower .zip.exe  [sha256: 37b894312d88c78d...]\n  C:\\Program Files (x86)\\Common Files\\microsoft shared\\indian fetish hot (!) vagina pregnant  (Tatjana).mpeg.exe  [sha256: 9d1dfad75b58aaa6...]\n  C:\\Program Files (x86)\\dotnet\\shared\\canadian gang bang licking .mpeg.exe  [sha256: baff50cc0ea5b167...]\n  C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\beast cumshot sleeping balls .rar.exe  [sha256: 6aec4abee30f53c1...]\n  C:\\Program Files (x86)\\Google\\Update\\Download\\fucking nude [free] mistress .zip.exe  [sha256: d089800b4b2ab998...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\tyrkish sperm kicking [milf] ejaculation .rar.exe  [sha256: e8d1b42dfe5d7386...]\n  C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\animal kicking lesbian .mpeg.exe  [sha256: 5cda67fbf339ae61...]\n  C:\\Program Files (x86)\\OpenOffice 4\\share\\russian bukkake [milf] feet .avi.exe  [sha256: b36a61bd055fa285...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\swedish nude [free] granny .zip.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\brasilian horse gay voyeur stockings  (Ashley).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\cum sleeping cock  (Britney).mpeg.exe\n  C:\\Documents and Settings\\All Users\\Templates\\xxx hardcore hidden  (Liz).mpeg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\canadian fucking public .rar.exe\n  C:\\Documents and Settings\\Default User\\Templates\\black horse [bangbus] young  (Melissa,Gina).mpeg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\canadian beastiality licking titts  (Sarah).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\canadian beastiality animal big hole bondage .zip.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\brasilian bukkake big stockings  (Jenna).avi.exe\n  C:\\Program Files\\Movie Maker\\Shared\\trambling lingerie [bangbus] cock sweet  (Christine).avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\tyrkish handjob lesbian boobs black hairunshaved  (Sylvia,Janette).zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\brasilian cumshot fucking voyeur femdom .mpeg.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\indian xxx sleeping girly .avi.exe\n  C:\\WINDOWS\\Temp\\swedish bukkake masturbation high heels .zip.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\gang bang hidden granny .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\african porn hidden hairy .mpg.exe\n  C:\\WINDOWS\\assembly\\temp\\brasilian xxx several models nipples mistress .avi.exe\n  C:\\WINDOWS\\assembly\\tmp\\gay voyeur ash  (Curtney,Sonja).mpg.exe\n  C:\\WINDOWS\\ime\\shared\\gay voyeur .mpg.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.65\n  [overlay]    entropy=7.880270957946777  size=895562\n",
  "ground_truth": {
    "hash": "41fd75ae1ced9aa2f43fb88c97c1da0819b1a10ac9f16fe1e5df0f7efee3e792",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9429,
    "malicious_count": 66,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "contain loop",
      "get disk information",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "send data on socket",
      "receive data on socket",
      "receive data",
      "send data",
      "initialize Winsock library",
      "resolve DNS",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\italian hardcore big .avi.exe",
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\WebResources\\Resource0\\static\\js\\plugins\\unified-share\\canadian beast porn public vagina shower .zip.exe",
        "C:\\Program Files (x86)\\Common Files\\microsoft shared\\indian fetish hot (!) vagina pregnant  (Tatjana).mpeg.exe",
        "C:\\Program Files (x86)\\Google\\Update\\1.3.34.11\\Recovery\\GUR7791.tmp\\beast cumshot sleeping balls .rar.exe",
        "C:\\Program Files (x86)\\Google\\Update\\Download\\fucking nude [free] mistress .zip.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\temp\\tyrkish sperm kicking [milf] ejaculation .rar.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\presets\\template\\animal kicking lesbian .mpeg.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\basic\\Template\\trambling full movie boobs .avi.exe",
        "C:\\Program Files (x86)\\OpenOffice 4\\share\\russian bukkake [milf] feet .avi.exe",
        "C:\\Program Files (x86)\\dotnet\\shared\\canadian gang bang licking .mpeg.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\xxx licking wifey .mpg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\french nude uncut legs hairy  (Ashley).mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\asian handjob licking gorgeoushorny .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\italian lingerie [milf] mistress .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "action sleeping black hairunshaved .mpeg.exe",
        "african porn hidden hairy .mpg.exe",
        "american lesbian porn hidden latex .mpeg.exe",
        "black bukkake big boobs .zip.exe",
        "black horse [bangbus] young  (Melissa,Gina).mpeg.exe",
        "brasilian bukkake big stockings  (Jenna).avi.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}