# dataset/preprocessed/423fe3492d006ed2f52fa537703077fd2bc1c9b85329cd283955d115b774790b.json
{
  "hash": "423fe3492d006ed2f52fa537703077fd2bc1c9b85329cd283955d115b774790b",
  "report_text": "--- NETWORK ACTIVITY ---\n\nDNS LOOKUPS (2):\n  VBOXSVR.ovh.net\n  vboxsvr.ovh.net\n\nSERVICES OPENED (3):\n  Messenger\n  NWCWorkstation\n  WebClient\n\n--- EXECUTION CHAIN ---\n\nPROCESS TREE:\n  american beastiality trambling catfight  (Samantha).mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\american beastiality trambling catfight  (Samantha).mpeg.exe\n    C:\\Users\\<USER>\\Downloads\\american beastiality trambling catfight  (Samantha).mpeg.exe\n\nCOMMAND EXECUTIONS (2):\n  ((null)) C:\\423fe3492d006ed2f52fa537703077fd2bc1c9b85329cd283955d115b774790b(null) [(null)]\n  C:\\Users\\<USER>\\Downloads\\american beastiality trambling catfight  (Samantha).mpeg.exe\n\n--- CAPABILITIES (CAPA) ---\n\nTECHNIQUE TRIGGERS (18):\n  - get disk information\n  - query or enumerate registry value\n  - persist via Run registry key\n  - get common file path\n  - get hostname\n  - enumerate network shares\n  - enumerate files on Windows\n  - enumerate files recursively\n  - start service\n  - modify service\n  - enumerate services\n  - stop service\n  - delete service\n  - enumerate gui resources\n  - enumerate processes\n  - link function at runtime on Windows\n  - accept command line arguments\n  - query environment variable\n\nCAPA RULE MATCHES (MBC):\n  - Data::Modulo [C0058]\n  - Operating System::Registry::Create Registry Key [C0036.004]\n  - Operating System::Registry::Open Registry Key [C0036.003]\n  - Anti-Behavioral Analysis::Dynamic Analysis Evasion::Delayed Execution [B0003.003]\n  - Process::Create Thread [C0038]\n  - Process::Check Mutex [C0043]\n  - File System::Create File [C0016]\n  - reference absolute stream path on Windows\n  - Process::Open Process [C0065]\n  - Process::Create Process [C0017]\n  - Communication::Socket Communication::Create UDP Socket [C0001.010]\n  - Memory::Allocate Memory [C0007]\n  - Discovery::System Information Discovery [E1082]\n  - contain loop\n  - Operating System::Registry::Query Registry Value [C0036.006]\n  - Persistence::Registry Run Keys / Startup Folder [F0012]\n  - Operating System::Registry::Set Registry Key [C0036.001]\n  - Discovery::File and Directory Discovery [E1083]\n  - connect network resource\n  - enumerate network shares\n  - Process::Create Mutex [C0042]\n  - get OS version\n  - File System::Read File [C0051]\n  - File System::Writes File [C0052]\n  - get service handle\n\n--- ARTIFACTS ---\n\nFILES DROPPED (30):\n  american kicking trambling public castration .mpg.exe  [sha256: 189992015b3e3285...]\n  black fetish beast hidden hole sm .avi.exe  [sha256: 29dc52877a003c92...]\n  blowjob uncut titts hairy  (Melissa).mpeg.exe  [sha256: 94e85d22c26e3359...]\n  brasilian beastiality blowjob sleeping femdom .avi.exe  [sha256: 16eb50f7a0b2f161...]\n  bukkake lesbian ejaculation .avi.exe  [sha256: 8263ae616dc2f77c...]\n  bukkake several models feet traffic .rar.exe  [sha256: 1335997401efe1d2...]\n  danish cumshot sperm girls feet .rar.exe  [sha256: e153f1fe9fbc4067...]\n  danish gang bang beast uncut .avi.exe  [sha256: ad274a655c60c922...]\n  driver.exe  [sha256: 675cce6f49339250...]\n  drwatson32.exe  [sha256: 13c7005b7e773616...]\n  fucking big shower .zip.exe  [sha256: ef1598bc1c7a758c...]\n  fucking full movie bedroom .rar.exe  [sha256: 97efea0ec905b097...]\n  fucking [milf] cock bedroom .zip.exe  [sha256: 8f63799ac97f6159...]\n  gay public .mpeg.exe  [sha256: 81a3591a818576c3...]\n  hardcore catfight balls .avi.exe  [sha256: 3a4c3c698a545059...]\n  hardcore hot (!) .mpeg.exe  [sha256: 8e91d7215d4f7734...]\n  horse hidden 40+ .zip.exe  [sha256: 7d0ac4223c3b678d...]\n  IE9setup.exe  [sha256: 39c8e0d860098c48...]\n  indian animal trambling hidden stockings .mpg.exe  [sha256: 419be62c99b7029a...]\n  indian gang bang bukkake several models bondage  (Sonja,Karin).zip.exe  [sha256: 224bc5c481b854b4...]\n\nFILES WRITTEN (30):\n  C:\\Documents and Settings\\<USER>\\Local Settings\\Temp\\russian handjob blowjob hidden .mpg.exe\n  C:\\Documents and Settings\\<USER>\\Templates\\tyrkish porn bukkake catfight .mpeg.exe\n  C:\\Documents and Settings\\All Users\\Application Data\\Microsoft\\Network\\Downloader\\tyrkish porn xxx catfight stockings .rar.exe\n  C:\\Documents and Settings\\All Users\\Templates\\american kicking trambling public castration .mpg.exe\n  C:\\Documents and Settings\\Default User\\Local Settings\\Temp\\russian fetish beast big circumcision .zip.exe\n  C:\\Documents and Settings\\Default User\\Templates\\indian nude lingerie [milf] feet bedroom  (Curtney).mpg.exe\n  C:\\Documents and Settings\\LocalService\\Local Settings\\Temp\\indian gang bang bukkake several models bondage  (Sonja,Karin).zip.exe\n  C:\\Documents and Settings\\NetworkService\\Local Settings\\Temp\\bukkake lesbian ejaculation .avi.exe\n  C:\\Program Files\\Common Files\\Microsoft Shared\\horse hidden 40+ .zip.exe\n  C:\\Program Files\\Movie Maker\\Shared\\black fetish beast hidden hole sm .avi.exe\n  C:\\WINDOWS\\Downloaded Program Files\\russian horse hardcore full movie .zip.exe\n  C:\\WINDOWS\\Microsoft.NET\\Framework\\v2.0.50727\\Temporary ASP.NET Files\\fucking big shower .zip.exe\n  C:\\WINDOWS\\SoftwareDistribution\\Download\\brasilian beastiality blowjob sleeping femdom .avi.exe\n  C:\\WINDOWS\\Temp\\japanese beastiality trambling [free] hotel  (Sonja,Karin).mpg.exe\n  C:\\WINDOWS\\WinSxS\\InstallTemp\\indian animal trambling hidden stockings .mpg.exe\n  C:\\WINDOWS\\assembly\\NativeImages_v2.0.50727_32\\Temp\\danish gang bang beast uncut .avi.exe\n  C:\\WINDOWS\\assembly\\temp\\italian handjob sperm [milf] .mpeg.exe\n  C:\\WINDOWS\\assembly\\tmp\\gay public .mpeg.exe\n  C:\\WINDOWS\\ime\\shared\\fucking full movie bedroom .rar.exe\n  C:\\WINDOWS\\mssrv.exe\n\nPERSISTENCE REGISTRY KEYS (1):\n  HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32 = C:\\Windows\\mssrv.exe\n\nMUTEXES CREATED (1):\n  mutex666\n\n--- FILE IDENTITY ---\n\nPE IMPORTS (DANGEROUS):\n  KERNEL32.dll: GetProcAddress, VirtualAlloc\n\nCOMPILER AND PACKER:\n  [Linker] Polink\n\nPE SECTIONS (entropy):\n  .text        entropy=6.45\n  .rdata       entropy=5.17\n  .data        entropy=7.64\n  [overlay]    entropy=7.8778839111328125  size=140701\n",
  "ground_truth": {
    "hash": "423fe3492d006ed2f52fa537703077fd2bc1c9b85329cd283955d115b774790b",
    "verdict": "MALICIOUS",
    "detection_rate": 0.9583,
    "malicious_count": 69,
    "popular_threat_category": [],
    "popular_threat_name": [],
    "suggested_label": null,
    "av_family_vote": "win32",
    "has_family_label": false,
    "ttps": [
      {
        "id": "T1007",
        "signature_description": "enumerate services"
      },
      {
        "id": "T1010",
        "signature_description": "enumerate gui resources"
      },
      {
        "id": "T1012",
        "signature_description": "query or enumerate registry value"
      },
      {
        "id": "T1057",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1059",
        "signature_description": "accept command line arguments"
      },
      {
        "id": "T1082",
        "signature_description": "get disk information"
      },
      {
        "id": "T1083",
        "signature_description": "get common file path"
      },
      {
        "id": "T1129",
        "signature_description": "link function at runtime on Windows"
      },
      {
        "id": "T1135",
        "signature_description": "enumerate network shares"
      },
      {
        "id": "T1489",
        "signature_description": "stop service"
      },
      {
        "id": "T1518",
        "signature_description": "enumerate processes"
      },
      {
        "id": "T1543",
        "signature_description": "start service"
      },
      {
        "id": "T1547",
        "signature_description": "persist via Run registry key"
      },
      {
        "id": "T1569",
        "signature_description": "modify service"
      }
    ],
    "capa_sig_names": [
      "calculate modulo 256 via x86 assembly",
      "create or open registry key",
      "delay execution",
      "create thread",
      "check mutex",
      "create or open file",
      "reference absolute stream path on Windows",
      "open process",
      "create process on Windows",
      "create UDP socket",
      "allocate memory",
      "allocate RW memory",
      "get disk information",
      "contain loop",
      "query or enumerate registry value",
      "persist via Run registry key",
      "set registry value",
      "get common file path",
      "connect network resource",
      "get hostname",
      "enumerate network shares",
      "create mutex",
      "get OS version",
      "read file via mapping",
      "write file on Windows",
      "enumerate files on Windows",
      "enumerate files recursively",
      "get service handle",
      "start service",
      "modify service",
      "enumerate services",
      "stop service",
      "delete service",
      "terminate process",
      "get graphical window text",
      "enumerate gui resources",
      "enumerate processes",
      "receive data on socket",
      "send data on socket",
      "receive data",
      "send data",
      "resolve DNS",
      "initialize Winsock library",
      "link function at runtime on Windows",
      "accept command line arguments",
      "query environment variable",
      "delete file",
      "read file on Windows",
      "reference analysis tools strings"
    ],
    "iocs": {
      "ips": [],
      "domains": [
        "vboxsvr.ovh.net"
      ],
      "urls": [],
      "registry_persistence": [
        "HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\mssrv32"
      ],
      "registry_all_count": 1,
      "files_dropped": [
        "C:\\Program Files (x86)\\Adobe\\Acrobat Reader DC\\Reader\\IDTemplates\\danish lesbian gay voyeur  (Melissa).zip.exe",
        "C:\\Program Files\\Common Files\\Microsoft Shared\\chinese lesbian blowjob lesbian shower .avi.exe",
        "C:\\Program Files\\DVD Maker\\Shared\\asian trambling girls .mpeg.exe",
        "C:\\Program Files\\Windows Journal\\Templates\\danish porn hidden .mpeg.exe",
        "C:\\Program Files\\Windows Sidebar\\Shared Gadgets\\french fetish sperm hot (!) granny .mpeg.exe",
        "C:\\Program Files\\dotnet\\shared\\gang bang nude [bangbus] vagina .avi.exe",
        "C:\\Windows\\mssrv.exe",
        "IE9setup.exe",
        "american kicking trambling public castration .mpg.exe",
        "black fetish beast hidden hole sm .avi.exe",
        "blowjob uncut titts hairy  (Melissa).mpeg.exe",
        "brasilian beastiality blowjob sleeping femdom .avi.exe",
        "bukkake lesbian ejaculation .avi.exe",
        "bukkake several models feet traffic .rar.exe",
        "danish cumshot sperm girls feet .rar.exe",
        "danish gang bang beast uncut .avi.exe",
        "driver.exe",
        "drwatson32.exe",
        "fucking [milf] cock bedroom .zip.exe",
        "fucking big shower .zip.exe",
        "fucking full movie bedroom .rar.exe",
        "gay public .mpeg.exe",
        "hardcore catfight balls .avi.exe",
        "hardcore hot (!) .mpeg.exe",
        "horse hidden 40+ .zip.exe",
        "indian animal trambling hidden stockings .mpg.exe",
        "indian gang bang bukkake several models bondage  (Sonja,Karin).zip.exe",
        "indian nude lingerie [milf] feet bedroom  (Curtney).mpg.exe",
        "indian porn trambling catfight hole 50+ .zip.exe",
        "italian action fucking several models .rar.exe",
        "italian handjob sperm [milf] .mpeg.exe",
        "japanese beastiality trambling [free] hotel  (Sonja,Karin).mpg.exe",
        "japanese horse horse catfight feet bondage .rar.exe",
        "mssrv.exe",
        "readme.txt.exe",
        "russian fetish beast big circumcision .zip.exe",
        "russian handjob blowjob hidden .mpg.exe",
        "russian horse hardcore full movie .zip.exe",
        "setup.exe",
        "start.exe",
        "tyrkish porn bukkake catfight .mpeg.exe",
        "tyrkish porn xxx catfight stockings .rar.exe",
        "undelete.exe",
        "xxx public titts leather .zip.exe"
      ],
      "mutexes": [
        "mutex666"
      ],
      "processes": []
    }
  },
  "num_sandboxes": 3,
  "sandbox_names": [
    "CAPA",
    "VirusTotal Cuckoofork",
    "VirusTotal Jujubox"
  ]
}